# extracted/clevis-21-1/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/clevis-luks-udisks2.desktop
usr/
usr/bin/
usr/bin/clevis
usr/bin/clevis-decrypt
usr/bin/clevis-decrypt-null
usr/bin/clevis-decrypt-sss
usr/bin/clevis-decrypt-tang
usr/bin/clevis-decrypt-tpm2
usr/bin/clevis-encrypt-null
usr/bin/clevis-encrypt-sss
usr/bin/clevis-encrypt-tang
usr/bin/clevis-encrypt-tpm2
usr/bin/clevis-luks-bind
usr/bin/clevis-luks-common-functions
usr/bin/clevis-luks-edit
usr/bin/clevis-luks-list
usr/bin/clevis-luks-pass
usr/bin/clevis-luks-regen
usr/bin/clevis-luks-report
usr/bin/clevis-luks-unbind
usr/bin/clevis-luks-unlock
usr/lib/
usr/lib/clevis-luks-askpass
usr/lib/clevis-luks-pkcs11-askpass
usr/lib/clevis-luks-pkcs11-askpin
usr/lib/clevis-luks-udisks2
usr/lib/clevis-luks-unlocker
usr/lib/dracut/
usr/lib/dracut/modules.d/
usr/lib/dracut/modules.d/60clevis-pin-null/
usr/lib/dracut/modules.d/60clevis-pin-null/module-setup.sh
usr/lib/dracut/modules.d/60clevis-pin-pkcs11/
usr/lib/dracut/modules.d/60clevis-pin-pkcs11/clevis-pkcs11-hook.sh
usr/lib/dracut/modules.d/60clevis-pin-pkcs11/clevis-pkcs11-prehook.sh
usr/lib/dracut/modules.d/60clevis-pin-pkcs11/module-setup.sh
usr/lib/dracut/modules.d/60clevis-pin-sss/
usr/lib/dracut/modules.d/60clevis-pin-sss/module-setup.sh
usr/lib/dracut/modules.d/60clevis-pin-tang/
usr/lib/dracut/modules.d/60clevis-pin-tang/module-setup.sh
usr/lib/dracut/modules.d/60clevis-pin-tpm2/
usr/lib/dracut/modules.d/60clevis-pin-tpm2/module-setup.sh
usr/lib/dracut/modules.d/60clevis/
usr/lib/dracut/modules.d/60clevis/clevis-hook.sh
usr/lib/dracut/modules.d/60clevis/module-setup.sh
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/clevis-luks-askpass.path
usr/lib/systemd/system/clevis-luks-askpass.service
usr/lib/systemd/system/clevis-luks-pkcs11-askpass.service
usr/lib/systemd/system/clevis-luks-pkcs11-askpass.socket
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/clevis
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/clevis-decrypt.1.gz
usr/share/man/man1/clevis-encrypt-sss.1.gz
usr/share/man/man1/clevis-encrypt-tang.1.gz
usr/share/man/man1/clevis-encrypt-tpm2.1.gz
usr/share/man/man1/clevis-luks-bind.1.gz
usr/share/man/man1/clevis-luks-edit.1.gz
usr/share/man/man1/clevis-luks-list.1.gz
usr/share/man/man1/clevis-luks-pass.1.gz
usr/share/man/man1/clevis-luks-regen.1.gz
usr/share/man/man1/clevis-luks-report.1.gz
usr/share/man/man1/clevis-luks-unbind.1.gz
usr/share/man/man1/clevis-luks-unlock.1.gz
usr/share/man/man1/clevis.1.gz
usr/share/man/man7/
usr/share/man/man7/clevis-luks-unlockers.7.gz
