# extracted/ddclient-3.11.2-5/files
%FILES%
etc/
etc/ddclient/
etc/ddclient/ddclient.conf
usr/
usr/bin/
usr/bin/ddclient
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ddclient.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/ddclient.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/ddclient.conf
usr/share/
usr/share/doc/
usr/share/doc/ddclient/
usr/share/doc/ddclient/README.cisco
usr/share/doc/ddclient/README.md
usr/share/licenses/
usr/share/licenses/ddclient/
usr/share/licenses/ddclient/COPYING
usr/share/licenses/ddclient/COPYRIGHT
var/
var/cache/
var/cache/ddclient/
