# extracted/woodpecker-server-2.7.1-1/files
%FILES%
etc/
etc/woodpecker/
etc/woodpecker/server.env
usr/
usr/bin/
usr/bin/woodpecker-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/woodpecker-server.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/woodpecker-server.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/woodpecker-server.conf
