# extracted/cronie-1.7.2-1/files
%FILES%
etc/
etc/anacrontab
etc/cron.d/
etc/cron.d/0hourly
etc/cron.daily/
etc/cron.deny
etc/cron.hourly/
etc/cron.hourly/0anacron
etc/cron.monthly/
etc/cron.weekly/
etc/crontab
etc/default/
etc/default/anacron
etc/pam.d/
etc/pam.d/crond
etc/sysconfig/
etc/sysconfig/crond
usr/
usr/bin/
usr/bin/anacron
usr/bin/crond
usr/bin/cronnext
usr/bin/crontab
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cronie.service
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/80-cronie.hook
usr/share/licenses/
usr/share/licenses/cronie/
usr/share/licenses/cronie/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cronnext.1.gz
usr/share/man/man1/crontab.1.gz
usr/share/man/man5/
usr/share/man/man5/anacrontab.5.gz
usr/share/man/man5/crontab.5.gz
usr/share/man/man8/
usr/share/man/man8/anacron.8.gz
usr/share/man/man8/cron.8.gz
usr/share/man/man8/crond.8.gz
var/
var/spool/
var/spool/anacron/
var/spool/cron/
