# extracted/prometheus-bird-exporter-1.4.3-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-bird-exporter
usr/
usr/bin/
usr/bin/prometheus-bird-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-bird-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-bird-exporter.conf
usr/share/
usr/share/doc/
usr/share/doc/prometheus-bird-exporter/
usr/share/doc/prometheus-bird-exporter/Dockerfile
usr/share/doc/prometheus-bird-exporter/README.md
usr/share/doc/prometheus-bird-exporter/daemonset.yaml
usr/share/licenses/
usr/share/licenses/prometheus-bird-exporter/
usr/share/licenses/prometheus-bird-exporter/LICENSE
