# extracted/dcd-1:0.15.2-2/files
%FILES%
etc/
etc/dcd.conf
usr/
usr/bin/
usr/bin/dcd-client
usr/bin/dcd-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dcd.service
usr/share/
usr/share/licenses/
usr/share/licenses/dcd/
usr/share/licenses/dcd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dcd-client.1.gz
usr/share/man/man1/dcd-server.1.gz
