# extracted/prometheus-mysqld-exporter-0.16.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-mysqld-exporter
usr/
usr/bin/
usr/bin/prometheus-mysqld-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-mysqld-exporter.service
