# extracted/boinc-8.0.2-2/files
%FILES%
usr/
usr/bin/
usr/bin/boinc
usr/bin/boinc_client
usr/bin/boinccmd
usr/bin/boincmgr
usr/include/
usr/include/boinc/
usr/include/boinc/app_ipc.h
usr/include/boinc/average.h
usr/include/boinc/base64.h
usr/include/boinc/boinc_api.h
usr/include/boinc/boinc_opencl.h
usr/include/boinc/boinc_stdio.h
usr/include/boinc/cal_boinc.h
usr/include/boinc/cc_config.h
usr/include/boinc/cert_sig.h
usr/include/boinc/cl_boinc.h
usr/include/boinc/common_defs.h
usr/include/boinc/coproc.h
usr/include/boinc/crypt.h
usr/include/boinc/diagnostics.h
usr/include/boinc/error_numbers.h
usr/include/boinc/filesys.h
usr/include/boinc/graphics2.h
usr/include/boinc/gui_rpc_client.h
usr/include/boinc/gutil.h
usr/include/boinc/hostinfo.h
usr/include/boinc/md5.h
usr/include/boinc/md5_file.h
usr/include/boinc/mem_usage.h
usr/include/boinc/mfile.h
usr/include/boinc/miofile.h
usr/include/boinc/msg_log.h
usr/include/boinc/network.h
usr/include/boinc/notice.h
usr/include/boinc/opencl_boinc.h
usr/include/boinc/parse.h
usr/include/boinc/prefs.h
usr/include/boinc/proc_control.h
usr/include/boinc/procinfo.h
usr/include/boinc/project_specific_defines.h
usr/include/boinc/proxy_info.h
usr/include/boinc/sched_msgs.h
usr/include/boinc/stackwalker_imports.h
usr/include/boinc/str_replace.h
usr/include/boinc/str_util.h
usr/include/boinc/svn_version.h
usr/include/boinc/url.h
usr/include/boinc/util.h
usr/include/boinc/version.h
usr/lib/
usr/lib/libboinc.a
usr/lib/libboinc_api.so
usr/lib/libboinc_api.so.8
usr/lib/libboinc_api.so.8.0.2
usr/lib/libboinc_crypt.a
usr/lib/libboinc_graphics2.so
usr/lib/libboinc_graphics2.so.8
usr/lib/libboinc_graphics2.so.8.0.2
usr/lib/libboinc_opencl.so
usr/lib/libboinc_opencl.so.8
usr/lib/libboinc_opencl.so.8.0.2
usr/lib/pkgconfig/
usr/lib/pkgconfig/libboinc.pc
usr/lib/pkgconfig/libboinc_api.pc
usr/lib/pkgconfig/libboinc_crypt.pc
usr/lib/pkgconfig/libboinc_graphics2.pc
usr/lib/pkgconfig/libboinc_opencl.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/boinc-client.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/boinc.conf
usr/share/
usr/share/applications/
usr/share/applications/boinc.desktop
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/boinc
usr/share/boinc-manager/
usr/share/boinc-manager/skins/
usr/share/boinc-manager/skins/Charity Engine/
usr/share/boinc-manager/skins/Charity Engine/background_image.png
usr/share/boinc-manager/skins/Charity Engine/ce_about.ico
usr/share/boinc-manager/skins/Charity Engine/ce_icon_play.png
usr/share/boinc-manager/skins/Charity Engine/ce_pause.png
usr/share/boinc-manager/skins/Charity Engine/ce_play.png
usr/share/boinc-manager/skins/Charity Engine/ce_stop.png
usr/share/boinc-manager/skins/Charity Engine/dialog_background_image.png
usr/share/boinc-manager/skins/Charity Engine/project_image.png
usr/share/boinc-manager/skins/Charity Engine/skin.xml
usr/share/boinc-manager/skins/Charity Engine/workunit_running_image.png
usr/share/boinc-manager/skins/Charity Engine/workunit_suspended_image.png
usr/share/boinc-manager/skins/Charity Engine/workunit_waiting_image.png
usr/share/boinc-manager/skins/Default/
usr/share/boinc-manager/skins/Default/background_image.png
usr/share/boinc-manager/skins/Default/skin.xml
usr/share/boinc-manager/skins/Default/workunit_running_image.png
usr/share/boinc-manager/skins/Default/workunit_suspended_image.png
usr/share/boinc-manager/skins/Default/workunit_waiting_image.png
usr/share/boinc-manager/skins/GridRepublic/
usr/share/boinc-manager/skins/GridRepublic/background_image.png
usr/share/boinc-manager/skins/GridRepublic/dialog_background_image.png
usr/share/boinc-manager/skins/GridRepublic/gr_about.ico
usr/share/boinc-manager/skins/GridRepublic/gr_icon_play.png
usr/share/boinc-manager/skins/GridRepublic/gr_pause.png
usr/share/boinc-manager/skins/GridRepublic/gr_play.png
usr/share/boinc-manager/skins/GridRepublic/gr_stop.png
usr/share/boinc-manager/skins/GridRepublic/project_image.png
usr/share/boinc-manager/skins/GridRepublic/skin.xml
usr/share/boinc-manager/skins/GridRepublic/workunit_running_image.png
usr/share/boinc-manager/skins/GridRepublic/workunit_suspended_image.png
usr/share/boinc-manager/skins/GridRepublic/workunit_waiting_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/advanced_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/attach_project_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/cancel_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/close_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/connecting_indicator_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/copy_all_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/copy_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/dialog_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/error_indicator_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/gauge_bg.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/help_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/left_arrow_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/messages_alert_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/messages_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/preferences_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/project_area_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/project_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/resume_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/right_arrow_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/save_button.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/spacer_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/state_indicator_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/suspend_link_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/tabArea_bg.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/wcg_about.ico
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/wcg_pause.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/wcg_play.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/wcg_stop.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_active_tab.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_animation_background_image copy.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_animation_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_area_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_gauge_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_gauge_progress_indicator_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_suspended_tab.png
usr/share/boinc-manager/skins/People for a Smarter Planet/graphic/workunit_tab_area_background_image.png
usr/share/boinc-manager/skins/People for a Smarter Planet/skin.xml
usr/share/boinc-manager/skins/ProgressThruProcessors/
usr/share/boinc-manager/skins/ProgressThruProcessors/background_image.png
usr/share/boinc-manager/skins/ProgressThruProcessors/dialog_background_image.png
usr/share/boinc-manager/skins/ProgressThruProcessors/project_image.png
usr/share/boinc-manager/skins/ProgressThruProcessors/ptp_about.ico
usr/share/boinc-manager/skins/ProgressThruProcessors/ptp_icon_play.png
usr/share/boinc-manager/skins/ProgressThruProcessors/ptp_pause.png
usr/share/boinc-manager/skins/ProgressThruProcessors/ptp_play.png
usr/share/boinc-manager/skins/ProgressThruProcessors/ptp_stop.png
usr/share/boinc-manager/skins/ProgressThruProcessors/skin.xml
usr/share/boinc-manager/skins/ProgressThruProcessors/workunit_running_image.png
usr/share/boinc-manager/skins/ProgressThruProcessors/workunit_suspended_image.png
usr/share/boinc-manager/skins/ProgressThruProcessors/workunit_waiting_image.png
usr/share/boinc-manager/skins/World Community Grid/
usr/share/boinc-manager/skins/World Community Grid/Green_dot.png
usr/share/boinc-manager/skins/World Community Grid/Red_dot.png
usr/share/boinc-manager/skins/World Community Grid/Yellow_dot.png
usr/share/boinc-manager/skins/World Community Grid/background_image.png
usr/share/boinc-manager/skins/World Community Grid/skin.xml
usr/share/boinc-manager/skins/World Community Grid/wcg.ico
usr/share/boinc-manager/skins/World Community Grid/wcg_32.png
usr/share/boinc-manager/skins/World Community Grid/wcg_50.png
usr/share/boinc-manager/skins/World Community Grid/wcg_about.ico
usr/share/boinc-manager/skins/World Community Grid/wcg_pause.png
usr/share/boinc-manager/skins/World Community Grid/wcg_play.png
usr/share/boinc-manager/skins/World Community Grid/wcg_stop.png
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/apps/
usr/share/icons/hicolor/64x64/apps/boinc.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/boinc.svg
usr/share/locale/
usr/share/locale/ar/
usr/share/locale/ar/LC_MESSAGES/
usr/share/locale/ar/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/be/
usr/share/locale/be/LC_MESSAGES/
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/bg/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ca/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/cs/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/da/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/de/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/es/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/fa/
usr/share/locale/fa/LC_MESSAGES/
usr/share/locale/fa_IR/
usr/share/locale/fa_IR/LC_MESSAGES/
usr/share/locale/fa_IR/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/fi/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/fr/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/he/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/hu/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/it_IT/
usr/share/locale/it_IT/LC_MESSAGES/
usr/share/locale/it_IT/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/it_IT/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ja/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ka/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ko/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/lv/
usr/share/locale/lv/LC_MESSAGES/
usr/share/locale/lv/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/lv/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ms/
usr/share/locale/ms/LC_MESSAGES/
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/nb/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/nl/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/nn/
usr/share/locale/nn/LC_MESSAGES/
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/pl/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/pt_BR/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/pt_PT/
usr/share/locale/pt_PT/LC_MESSAGES/
usr/share/locale/pt_PT/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/pt_PT/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ro/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/ru/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/sk/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/sv/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/tr/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/uk/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/zh_CN/LC_MESSAGES/BOINC-Manager.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/BOINC-Client.mo
usr/share/locale/zh_TW/LC_MESSAGES/BOINC-Manager.mo
usr/share/pixmaps/
usr/share/pixmaps/boinc.png
