# extracted/caddy-2.8.4-1/files
%FILES%
etc/
etc/caddy/
etc/caddy/Caddyfile
etc/caddy/conf.d/
usr/
usr/bin/
usr/bin/caddy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/caddy-api.service
usr/lib/systemd/system/caddy.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/caddy.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/caddy.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/caddy
usr/share/caddy/
usr/share/caddy/index.html
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/caddy.fish
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_caddy
