# extracted/sonic-1.4.0-1/files
%FILES%
etc/
etc/sonic.cfg
usr/
usr/bin/
usr/bin/sonic
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sonic.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/sonic.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sonic.conf
usr/share/
usr/share/doc/
usr/share/doc/sonic/
usr/share/doc/sonic/CHANGELOG.md
usr/share/doc/sonic/CODE_OF_CONDUCT.md
usr/share/doc/sonic/CONFIGURATION.md
usr/share/doc/sonic/CONTRIBUTING.md
usr/share/doc/sonic/INNER_WORKINGS.md
usr/share/doc/sonic/LICENSE.md
usr/share/doc/sonic/PACKAGING.md
usr/share/doc/sonic/PROTOCOL.md
usr/share/doc/sonic/README.md
