# extracted/prosody-1:0.12.4-2/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/prosody
etc/prosody/
etc/prosody/certs/
etc/prosody/migrator.cfg.lua
etc/prosody/prosody.cfg.lua
usr/
usr/bin/
usr/bin/prosody
usr/bin/prosody-migrator
usr/bin/prosodyctl
usr/lib/
usr/lib/prosody/
usr/lib/prosody/core/
usr/lib/prosody/core/certmanager.lua
usr/lib/prosody/core/configmanager.lua
usr/lib/prosody/core/features.lua
usr/lib/prosody/core/hostmanager.lua
usr/lib/prosody/core/loggingmanager.lua
usr/lib/prosody/core/moduleapi.lua
usr/lib/prosody/core/modulemanager.lua
usr/lib/prosody/core/portmanager.lua
usr/lib/prosody/core/rostermanager.lua
usr/lib/prosody/core/s2smanager.lua
usr/lib/prosody/core/sessionmanager.lua
usr/lib/prosody/core/stanza_router.lua
usr/lib/prosody/core/statsmanager.lua
usr/lib/prosody/core/storagemanager.lua
usr/lib/prosody/core/usermanager.lua
usr/lib/prosody/modules/
usr/lib/prosody/modules/adhoc/
usr/lib/prosody/modules/adhoc/adhoc.lib.lua
usr/lib/prosody/modules/adhoc/mod_adhoc.lua
usr/lib/prosody/modules/mod_admin_adhoc.lua
usr/lib/prosody/modules/mod_admin_shell.lua
usr/lib/prosody/modules/mod_admin_socket.lua
usr/lib/prosody/modules/mod_admin_telnet.lua
usr/lib/prosody/modules/mod_announce.lua
usr/lib/prosody/modules/mod_auth_anonymous.lua
usr/lib/prosody/modules/mod_auth_insecure.lua
usr/lib/prosody/modules/mod_auth_internal_hashed.lua
usr/lib/prosody/modules/mod_auth_internal_plain.lua
usr/lib/prosody/modules/mod_auth_ldap.lua
usr/lib/prosody/modules/mod_authz_internal.lua
usr/lib/prosody/modules/mod_blocklist.lua
usr/lib/prosody/modules/mod_bookmarks.lua
usr/lib/prosody/modules/mod_bosh.lua
usr/lib/prosody/modules/mod_c2s.lua
usr/lib/prosody/modules/mod_carbons.lua
usr/lib/prosody/modules/mod_component.lua
usr/lib/prosody/modules/mod_cron.lua
usr/lib/prosody/modules/mod_csi.lua
usr/lib/prosody/modules/mod_csi_simple.lua
usr/lib/prosody/modules/mod_debug_sql.lua
usr/lib/prosody/modules/mod_dialback.lua
usr/lib/prosody/modules/mod_disco.lua
usr/lib/prosody/modules/mod_external_services.lua
usr/lib/prosody/modules/mod_groups.lua
usr/lib/prosody/modules/mod_http.lua
usr/lib/prosody/modules/mod_http_errors.lua
usr/lib/prosody/modules/mod_http_file_share.lua
usr/lib/prosody/modules/mod_http_files.lua
usr/lib/prosody/modules/mod_http_openmetrics.lua
usr/lib/prosody/modules/mod_invites.lua
usr/lib/prosody/modules/mod_invites_adhoc.lua
usr/lib/prosody/modules/mod_invites_register.lua
usr/lib/prosody/modules/mod_iq.lua
usr/lib/prosody/modules/mod_lastactivity.lua
usr/lib/prosody/modules/mod_legacyauth.lua
usr/lib/prosody/modules/mod_limits.lua
usr/lib/prosody/modules/mod_mam/
usr/lib/prosody/modules/mod_mam/mamprefs.lib.lua
usr/lib/prosody/modules/mod_mam/mamprefsxml.lib.lua
usr/lib/prosody/modules/mod_mam/mod_mam.lua
usr/lib/prosody/modules/mod_message.lua
usr/lib/prosody/modules/mod_mimicking.lua
usr/lib/prosody/modules/mod_motd.lua
usr/lib/prosody/modules/mod_muc_mam.lua
usr/lib/prosody/modules/mod_muc_unique.lua
usr/lib/prosody/modules/mod_net_multiplex.lua
usr/lib/prosody/modules/mod_offline.lua
usr/lib/prosody/modules/mod_pep.lua
usr/lib/prosody/modules/mod_pep_plus.lua
usr/lib/prosody/modules/mod_pep_simple.lua
usr/lib/prosody/modules/mod_ping.lua
usr/lib/prosody/modules/mod_posix.lua
usr/lib/prosody/modules/mod_presence.lua
usr/lib/prosody/modules/mod_private.lua
usr/lib/prosody/modules/mod_proxy65.lua
usr/lib/prosody/modules/mod_pubsub/
usr/lib/prosody/modules/mod_pubsub/mod_pubsub.lua
usr/lib/prosody/modules/mod_pubsub/pubsub.lib.lua
usr/lib/prosody/modules/mod_register.lua
usr/lib/prosody/modules/mod_register_ibr.lua
usr/lib/prosody/modules/mod_register_limits.lua
usr/lib/prosody/modules/mod_roster.lua
usr/lib/prosody/modules/mod_s2s.lua
usr/lib/prosody/modules/mod_s2s_auth_certs.lua
usr/lib/prosody/modules/mod_s2s_bidi.lua
usr/lib/prosody/modules/mod_saslauth.lua
usr/lib/prosody/modules/mod_scansion_record.lua
usr/lib/prosody/modules/mod_server_contact_info.lua
usr/lib/prosody/modules/mod_smacks.lua
usr/lib/prosody/modules/mod_stanza_debug.lua
usr/lib/prosody/modules/mod_storage_internal.lua
usr/lib/prosody/modules/mod_storage_memory.lua
usr/lib/prosody/modules/mod_storage_none.lua
usr/lib/prosody/modules/mod_storage_sql.lua
usr/lib/prosody/modules/mod_storage_xep0227.lua
usr/lib/prosody/modules/mod_time.lua
usr/lib/prosody/modules/mod_tls.lua
usr/lib/prosody/modules/mod_tokenauth.lua
usr/lib/prosody/modules/mod_tombstones.lua
usr/lib/prosody/modules/mod_turn_external.lua
usr/lib/prosody/modules/mod_unknown.lua
usr/lib/prosody/modules/mod_uptime.lua
usr/lib/prosody/modules/mod_user_account_management.lua
usr/lib/prosody/modules/mod_vcard.lua
usr/lib/prosody/modules/mod_vcard4.lua
usr/lib/prosody/modules/mod_vcard_legacy.lua
usr/lib/prosody/modules/mod_version.lua
usr/lib/prosody/modules/mod_watchregistrations.lua
usr/lib/prosody/modules/mod_websocket.lua
usr/lib/prosody/modules/mod_welcome.lua
usr/lib/prosody/modules/mod_windows.lua
usr/lib/prosody/modules/muc/
usr/lib/prosody/modules/muc/config_form_sections.lib.lua
usr/lib/prosody/modules/muc/description.lib.lua
usr/lib/prosody/modules/muc/hats.lib.lua
usr/lib/prosody/modules/muc/hidden.lib.lua
usr/lib/prosody/modules/muc/history.lib.lua
usr/lib/prosody/modules/muc/language.lib.lua
usr/lib/prosody/modules/muc/lock.lib.lua
usr/lib/prosody/modules/muc/members_only.lib.lua
usr/lib/prosody/modules/muc/mod_muc.lua
usr/lib/prosody/modules/muc/moderated.lib.lua
usr/lib/prosody/modules/muc/muc.lib.lua
usr/lib/prosody/modules/muc/name.lib.lua
usr/lib/prosody/modules/muc/occupant.lib.lua
usr/lib/prosody/modules/muc/occupant_id.lib.lua
usr/lib/prosody/modules/muc/password.lib.lua
usr/lib/prosody/modules/muc/persistent.lib.lua
usr/lib/prosody/modules/muc/presence_broadcast.lib.lua
usr/lib/prosody/modules/muc/register.lib.lua
usr/lib/prosody/modules/muc/request.lib.lua
usr/lib/prosody/modules/muc/subject.lib.lua
usr/lib/prosody/modules/muc/util.lib.lua
usr/lib/prosody/modules/muc/whois.lib.lua
usr/lib/prosody/net/
usr/lib/prosody/net/adns.lua
usr/lib/prosody/net/connect.lua
usr/lib/prosody/net/cqueues.lua
usr/lib/prosody/net/dns.lua
usr/lib/prosody/net/http.lua
usr/lib/prosody/net/http/
usr/lib/prosody/net/http/codes.lua
usr/lib/prosody/net/http/errors.lua
usr/lib/prosody/net/http/files.lua
usr/lib/prosody/net/http/parser.lua
usr/lib/prosody/net/http/server.lua
usr/lib/prosody/net/httpserver.lua
usr/lib/prosody/net/resolvers/
usr/lib/prosody/net/resolvers/basic.lua
usr/lib/prosody/net/resolvers/chain.lua
usr/lib/prosody/net/resolvers/manual.lua
usr/lib/prosody/net/resolvers/service.lua
usr/lib/prosody/net/server.lua
usr/lib/prosody/net/server_epoll.lua
usr/lib/prosody/net/server_event.lua
usr/lib/prosody/net/server_select.lua
usr/lib/prosody/net/stun.lua
usr/lib/prosody/net/unbound.lua
usr/lib/prosody/net/websocket.lua
usr/lib/prosody/net/websocket/
usr/lib/prosody/net/websocket/frames.lua
usr/lib/prosody/prosody.version
usr/lib/prosody/util/
usr/lib/prosody/util/adhoc.lua
usr/lib/prosody/util/adminstream.lua
usr/lib/prosody/util/argparse.lua
usr/lib/prosody/util/array.lua
usr/lib/prosody/util/async.lua
usr/lib/prosody/util/bit53.lua
usr/lib/prosody/util/bitcompat.lua
usr/lib/prosody/util/cache.lua
usr/lib/prosody/util/caps.lua
usr/lib/prosody/util/compat.so
usr/lib/prosody/util/dataforms.lua
usr/lib/prosody/util/datamanager.lua
usr/lib/prosody/util/datamapper.lua
usr/lib/prosody/util/datetime.lua
usr/lib/prosody/util/dbuffer.lua
usr/lib/prosody/util/debug.lua
usr/lib/prosody/util/dependencies.lua
usr/lib/prosody/util/dns.lua
usr/lib/prosody/util/dnsregistry.lua
usr/lib/prosody/util/encodings.so
usr/lib/prosody/util/envload.lua
usr/lib/prosody/util/error.lua
usr/lib/prosody/util/events.lua
usr/lib/prosody/util/filters.lua
usr/lib/prosody/util/format.lua
usr/lib/prosody/util/gc.lua
usr/lib/prosody/util/hashes.so
usr/lib/prosody/util/hashring.lua
usr/lib/prosody/util/helpers.lua
usr/lib/prosody/util/hex.lua
usr/lib/prosody/util/hmac.lua
usr/lib/prosody/util/http.lua
usr/lib/prosody/util/human/
usr/lib/prosody/util/human/io.lua
usr/lib/prosody/util/human/units.lua
usr/lib/prosody/util/id.lua
usr/lib/prosody/util/import.lua
usr/lib/prosody/util/indexedbheap.lua
usr/lib/prosody/util/interpolation.lua
usr/lib/prosody/util/ip.lua
usr/lib/prosody/util/iterators.lua
usr/lib/prosody/util/jid.lua
usr/lib/prosody/util/json.lua
usr/lib/prosody/util/jsonpointer.lua
usr/lib/prosody/util/jsonschema.lua
usr/lib/prosody/util/jwt.lua
usr/lib/prosody/util/logger.lua
usr/lib/prosody/util/mercurial.lua
usr/lib/prosody/util/multitable.lua
usr/lib/prosody/util/net.so
usr/lib/prosody/util/openmetrics.lua
usr/lib/prosody/util/openssl.lua
usr/lib/prosody/util/paths.lua
usr/lib/prosody/util/pluginloader.lua
usr/lib/prosody/util/poll.so
usr/lib/prosody/util/pposix.so
usr/lib/prosody/util/presence.lua
usr/lib/prosody/util/promise.lua
usr/lib/prosody/util/prosodyctl.lua
usr/lib/prosody/util/prosodyctl/
usr/lib/prosody/util/prosodyctl/cert.lua
usr/lib/prosody/util/prosodyctl/check.lua
usr/lib/prosody/util/prosodyctl/shell.lua
usr/lib/prosody/util/pubsub.lua
usr/lib/prosody/util/queue.lua
usr/lib/prosody/util/random.lua
usr/lib/prosody/util/rfc6724.lua
usr/lib/prosody/util/ringbuffer.so
usr/lib/prosody/util/rsm.lua
usr/lib/prosody/util/sasl.lua
usr/lib/prosody/util/sasl/
usr/lib/prosody/util/sasl/anonymous.lua
usr/lib/prosody/util/sasl/external.lua
usr/lib/prosody/util/sasl/plain.lua
usr/lib/prosody/util/sasl/scram.lua
usr/lib/prosody/util/serialization.lua
usr/lib/prosody/util/session.lua
usr/lib/prosody/util/set.lua
usr/lib/prosody/util/signal.so
usr/lib/prosody/util/smqueue.lua
usr/lib/prosody/util/sql.lua
usr/lib/prosody/util/sslconfig.lua
usr/lib/prosody/util/stanza.lua
usr/lib/prosody/util/startup.lua
usr/lib/prosody/util/statistics.lua
usr/lib/prosody/util/statsd.lua
usr/lib/prosody/util/strbitop.so
usr/lib/prosody/util/struct.so
usr/lib/prosody/util/table.so
usr/lib/prosody/util/template.lua
usr/lib/prosody/util/termcolours.lua
usr/lib/prosody/util/throttle.lua
usr/lib/prosody/util/time.so
usr/lib/prosody/util/timer.lua
usr/lib/prosody/util/uuid.lua
usr/lib/prosody/util/vcard.lua
usr/lib/prosody/util/watchdog.lua
usr/lib/prosody/util/x509.lua
usr/lib/prosody/util/xml.lua
usr/lib/prosody/util/xmppstream.lua
usr/lib/prosody/util/xpcall.lua
usr/lib/prosody/util/xtemplate.lua
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prosody.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prosody.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/prosody.conf
usr/share/
usr/share/licenses/
usr/share/licenses/prosody/
usr/share/licenses/prosody/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/prosodyctl.1.gz
usr/share/prosody/
usr/share/prosody/tools/
usr/share/prosody/tools/cfgdump.lua
usr/share/prosody/tools/dnsregistry.lua
usr/share/prosody/tools/ejabberd2prosody.lua
usr/share/prosody/tools/ejabberdsql2prosody.lua
usr/share/prosody/tools/erlparse.lua
usr/share/prosody/tools/form2table.lua
usr/share/prosody/tools/generate_format_spec.lua
usr/share/prosody/tools/http-status-codes.lua
usr/share/prosody/tools/jabberd14sql2prosody.lua
usr/share/prosody/tools/linedebug.lua
usr/share/prosody/tools/make_repo.lua
usr/share/prosody/tools/modtrace.lua
usr/share/prosody/tools/openfire2prosody.lua
var/
var/lib/
var/lib/prosody/
