# extracted/uwsgi-2.0.28-2/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/uwsgi
etc/uwsgi/
etc/uwsgi/emperor.ini
usr/
usr/bin/
usr/bin/uwsgi
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/emperor.uwsgi.service
usr/lib/systemd/system/emperor.uwsgi.socket
usr/lib/systemd/system/uwsgi@.service
usr/lib/systemd/system/uwsgi@.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/uwsgi.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/uwsgi.conf
usr/share/
usr/share/licenses/
usr/share/licenses/uwsgi/
usr/share/licenses/uwsgi/LICENSE
