# extracted/minio-2024.12.13-1/files
%FILES%
etc/
etc/minio/
etc/minio/minio.conf
srv/
srv/minio/
usr/
usr/bin/
usr/bin/minio
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/minio.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/minio.conf
var/
var/lib/
var/lib/minio/
