# extracted/prometheus-nut-exporter-3.1.1-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-nut-exporter
usr/
usr/bin/
usr/bin/prometheus-nut-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-nut-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-nut-exporter.conf
