# extracted/hefur-1.0-19/files
%FILES%
usr/
usr/bin/
usr/bin/hefurctl
usr/bin/hefurd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hefurd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/hefurd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/hefurd.conf
usr/share/
usr/share/doc/
usr/share/doc/hefur/
usr/share/doc/hefur/hefur/
usr/share/doc/hefur/hefur/manual.md
usr/share/hefur/
usr/share/hefur/www/
usr/share/hefur/www/css/
usr/share/hefur/www/css/bootstrap-grid.css
usr/share/hefur/www/css/bootstrap-grid.css.map
usr/share/hefur/www/css/bootstrap-grid.min.css
usr/share/hefur/www/css/bootstrap-grid.min.css.map
usr/share/hefur/www/css/bootstrap-reboot.css
usr/share/hefur/www/css/bootstrap-reboot.css.map
usr/share/hefur/www/css/bootstrap-reboot.min.css
usr/share/hefur/www/css/bootstrap-reboot.min.css.map
usr/share/hefur/www/css/bootstrap.css
usr/share/hefur/www/css/bootstrap.css.map
usr/share/hefur/www/css/bootstrap.min.css
usr/share/hefur/www/css/bootstrap.min.css.map
usr/share/hefur/www/img/
usr/share/hefur/www/img/favicon.ico
usr/share/hefur/www/img/magnet-icon.gif
usr/share/hefur/www/img/peers.png
usr/share/hefur/www/img/torrent-icon.png
usr/share/hefur/www/js/
usr/share/hefur/www/js/bootstrap.bundle.js
usr/share/hefur/www/js/bootstrap.bundle.js.map
usr/share/hefur/www/js/bootstrap.bundle.min.js
usr/share/hefur/www/js/bootstrap.bundle.min.js.map
usr/share/hefur/www/js/bootstrap.js
usr/share/hefur/www/js/bootstrap.js.map
usr/share/hefur/www/js/bootstrap.min.js
usr/share/hefur/www/js/bootstrap.min.js.map
usr/share/hefur/www/js/table-sorter.js
usr/share/hefur/www/tpl/
usr/share/hefur/www/tpl/page.html
usr/share/hefur/www/tpl/peers.html
usr/share/hefur/www/tpl/stat.html
usr/share/licenses/
usr/share/licenses/hefur/
usr/share/licenses/hefur/LICENSE.md
