# extracted/syncthing-1.28.1-1/files
%FILES%
etc/
etc/ufw/
etc/ufw/applications.d/
etc/ufw/applications.d/ufw-syncthing
usr/
usr/bin/
usr/bin/syncthing
usr/lib/
usr/lib/sysctl.d/
usr/lib/sysctl.d/30-syncthing.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/syncthing@.service
usr/lib/systemd/user/
usr/lib/systemd/user/syncthing.service
usr/share/
usr/share/doc/
usr/share/doc/syncthing/
usr/share/doc/syncthing/README.md
usr/share/licenses/
usr/share/licenses/syncthing/
usr/share/licenses/syncthing/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/syncthing.1.gz
usr/share/man/man5/
usr/share/man/man5/syncthing-config.5.gz
usr/share/man/man5/syncthing-stignore.5.gz
usr/share/man/man7/
usr/share/man/man7/syncthing-bep.7.gz
usr/share/man/man7/syncthing-device-ids.7.gz
usr/share/man/man7/syncthing-event-api.7.gz
usr/share/man/man7/syncthing-faq.7.gz
usr/share/man/man7/syncthing-globaldisco.7.gz
usr/share/man/man7/syncthing-localdisco.7.gz
usr/share/man/man7/syncthing-networking.7.gz
usr/share/man/man7/syncthing-relay.7.gz
usr/share/man/man7/syncthing-rest-api.7.gz
usr/share/man/man7/syncthing-security.7.gz
usr/share/man/man7/syncthing-versioning.7.gz
