# extracted/acpid-2.0.34-2/files
%FILES%
etc/
etc/acpi/
etc/acpi/events/
etc/acpi/events/anything
etc/acpi/handler.sh
usr/
usr/bin/
usr/bin/acpi_listen
usr/bin/acpid
usr/bin/kacpimon
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/acpid.service
usr/share/
usr/share/doc/
usr/share/doc/acpid/
usr/share/doc/acpid/COPYING
usr/share/doc/acpid/Changelog
usr/share/doc/acpid/README
usr/share/doc/acpid/TESTPLAN
usr/share/doc/acpid/TODO
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/acpi_listen.8.gz
usr/share/man/man8/acpid.8.gz
usr/share/man/man8/kacpimon.8.gz
