# extracted/ssh-tpm-agent-0.7.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/ssh-tpm-add
usr/bin/ssh-tpm-agent
usr/bin/ssh-tpm-hostkeys
usr/bin/ssh-tpm-keygen
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ssh-tpm-agent.service
usr/lib/systemd/system/ssh-tpm-agent.socket
usr/lib/systemd/system/ssh-tpm-genkeys.service
usr/lib/systemd/user/
usr/lib/systemd/user/ssh-tpm-agent.service
usr/lib/systemd/user/ssh-tpm-agent.socket
