# extracted/udp2raw-20230206.0-2/files
%FILES%
etc/
etc/udp2raw/
etc/udp2raw/example.conf
usr/
usr/bin/
usr/bin/udp2raw
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/udp2raw@.service
usr/lib/udp2raw/
usr/lib/udp2raw/udp2raw_script.sh
usr/share/
usr/share/licenses/
usr/share/licenses/udp2raw/
usr/share/licenses/udp2raw/LICENSE.md
