# extracted/cpupower-6.12-8/files
%FILES%
etc/
etc/cpufreq-bench.conf
etc/default/
etc/default/cpupower
usr/
usr/bin/
usr/bin/cpufreq-bench
usr/bin/cpufreq-bench_plot.sh
usr/bin/cpupower
usr/include/
usr/include/cpufreq.h
usr/include/cpuidle.h
usr/include/powercap.h
usr/lib/
usr/lib/libcpupower.so
usr/lib/libcpupower.so.0.0.1
usr/lib/libcpupower.so.1
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/cpupower
usr/lib/systemd/system/
usr/lib/systemd/system/cpupower.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/cpupower
usr/share/doc/
usr/share/doc/cpupower/
usr/share/doc/cpupower/README-BENCH
usr/share/doc/cpupower/cpufreq-bench_script.sh
usr/share/locale/
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/cpupower.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/cpupower.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/cpupower.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/cpupower.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/cpupower.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/cpupower.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cpupower-frequency-info.1.gz
usr/share/man/man1/cpupower-frequency-set.1.gz
usr/share/man/man1/cpupower-idle-info.1.gz
usr/share/man/man1/cpupower-idle-set.1.gz
usr/share/man/man1/cpupower-info.1.gz
usr/share/man/man1/cpupower-monitor.1.gz
usr/share/man/man1/cpupower-powercap-info.1.gz
usr/share/man/man1/cpupower-set.1.gz
usr/share/man/man1/cpupower.1.gz
