# extracted/pgbouncer-1.22.1-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/pgbouncer
etc/pgbouncer/
etc/pgbouncer/pgbouncer.ini
usr/
usr/bin/
usr/bin/pgbouncer
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pgbouncer.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pgbouncer.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pgbouncer.conf
usr/share/
usr/share/doc/
usr/share/doc/pgbouncer/
usr/share/doc/pgbouncer/NEWS.md
usr/share/doc/pgbouncer/README.md
usr/share/doc/pgbouncer/pgbouncer-minimal.ini
usr/share/doc/pgbouncer/pgbouncer.ini
usr/share/doc/pgbouncer/pgbouncer.service
usr/share/doc/pgbouncer/pgbouncer.socket
usr/share/doc/pgbouncer/userlist.txt
usr/share/licenses/
usr/share/licenses/pgbouncer/
usr/share/licenses/pgbouncer/COPYRIGHT
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/pgbouncer.1.gz
usr/share/man/man5/
usr/share/man/man5/pgbouncer.5.gz
