# extracted/meilisearch-1.11.2-1/files
%FILES%
etc/
etc/meilisearch.conf
usr/
usr/bin/
usr/bin/meilisearch
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/meilisearch.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/meilisearch.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/meilisearch.conf
usr/share/
usr/share/licenses/
usr/share/licenses/meilisearch/
usr/share/licenses/meilisearch/LICENSE
