# extracted/postfwd-2.03-2/files
%FILES%
etc/
etc/postfwd/
etc/postfwd/postfwd-old.cf.sample
etc/postfwd/postfwd.cf
etc/postfwd/postfwd.cf.sample
etc/postfwd/postfwd.settings.sample
etc/postfwd/postfwd.settings.sample.full
usr/
usr/bin/
usr/bin/postfwd
usr/bin/postfwd-script.sh
usr/bin/postfwd1
usr/bin/postfwd2
usr/bin/postfwd3
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/postfwd.service
usr/lib/systemd/system/postfwd2.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/postfwd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/postfwd/
usr/share/licenses/postfwd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/postfwd.8.gz
usr/share/man/man8/postfwd1.8.gz
usr/share/man/man8/postfwd2.8.gz
usr/share/man/man8/postfwd3.8.gz
usr/share/postfwd/
usr/share/postfwd/arch.html
usr/share/postfwd/docker.txt
usr/share/postfwd/hapolicy.html
usr/share/postfwd/hapolicy.txt
usr/share/postfwd/hapolicy/
usr/share/postfwd/hapolicy/hapolicy
usr/share/postfwd/hapolicy/hapolicy.html
usr/share/postfwd/hapolicy/hapolicy.txt
usr/share/postfwd/hapolicy/hapolicy01.png
usr/share/postfwd/hapolicy/hapolicy02.png
usr/share/postfwd/old/
usr/share/postfwd/old/postfwd1.CHANGELOG
usr/share/postfwd/old/postfwd2.CHANGELOG
usr/share/postfwd/postfwd-ARCH.png
usr/share/postfwd/postfwd-client.pl
usr/share/postfwd/postfwd.CHANGELOG
usr/share/postfwd/postfwd.html
usr/share/postfwd/postfwd.txt
usr/share/postfwd/postfwd3.CHANGELOG
usr/share/postfwd/postfwd3.html
usr/share/postfwd/postfwd3.txt
usr/share/postfwd/quick.html
usr/share/postfwd/rblcheck.pl
usr/share/postfwd/request.sample
usr/share/postfwd/versions.html
