# extracted/zabbix-proxy-7.2.0-1/files
%FILES%
etc/
etc/zabbix/
etc/zabbix/zabbix_proxy.conf
usr/
usr/bin/
usr/bin/zabbix_proxy_mysql
usr/bin/zabbix_proxy_postgresql
usr/bin/zabbix_proxy_sqlite3
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/zabbix-proxy-mysql.service
usr/lib/systemd/system/zabbix-proxy-pgsql.service
usr/lib/systemd/system/zabbix-proxy-sqlite.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/zabbix-proxy.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/zabbix-proxy.conf
usr/share/
usr/share/zabbix-proxy/
usr/share/zabbix-proxy/mysql/
usr/share/zabbix-proxy/mysql/schema.sql
usr/share/zabbix-proxy/postgresql/
usr/share/zabbix-proxy/postgresql/schema.sql
usr/share/zabbix-proxy/sqlite3/
usr/share/zabbix-proxy/sqlite3/schema.sql
