# extracted/lirc-1:0.10.2-4/files
%FILES%
etc/
etc/lirc/
etc/lirc/irexec.lircrc
etc/lirc/lirc_options.conf
etc/lirc/lircd.conf
etc/lirc/lircd.conf.d/
etc/lirc/lircd.conf.d/README.conf.d
etc/lirc/lircd.conf.d/devinput.lircd.conf
etc/lirc/lircmd.conf
etc/logrotate.d/
etc/logrotate.d/lirc
usr/
usr/bin/
usr/bin/ircat
usr/bin/irdb-get
usr/bin/irexec
usr/bin/irpipe
usr/bin/irpty
usr/bin/irrecord
usr/bin/irsend
usr/bin/irsimreceive
usr/bin/irsimsend
usr/bin/irtestcase
usr/bin/irtext2udp
usr/bin/irw
usr/bin/irxevent
usr/bin/lirc-config-tool
usr/bin/lirc-data2table
usr/bin/lirc-init-db
usr/bin/lirc-lsplugins
usr/bin/lirc-lsremotes
usr/bin/lirc-make-devinput
usr/bin/lirc-postinstall
usr/bin/lirc-setup
usr/bin/lircd
usr/bin/lircd-setup
usr/bin/lircd-uinput
usr/bin/lircmd
usr/bin/lircrcd
usr/bin/mode2
usr/bin/pronto2lirc
usr/bin/xmode2
usr/include/
usr/include/lirc/
usr/include/lirc/ciniparser.h
usr/include/lirc/config.h
usr/include/lirc/config_file.h
usr/include/lirc/config_flags.h
usr/include/lirc/curl_poll.h
usr/include/lirc/dictionary.h
usr/include/lirc/driver.h
usr/include/lirc/drv_admin.h
usr/include/lirc/drv_enum.h
usr/include/lirc/dump_config.h
usr/include/lirc/input-event-codes.h
usr/include/lirc/input_map.h
usr/include/lirc/input_map.inc
usr/include/lirc/ir_remote.h
usr/include/lirc/ir_remote_types.h
usr/include/lirc/irpipe.h
usr/include/lirc/irrecord.h
usr/include/lirc/line_buffer.h
usr/include/lirc/lirc-utils.h
usr/include/lirc/lirc_client.h
usr/include/lirc/lirc_config.h
usr/include/lirc/lirc_log.h
usr/include/lirc/lirc_options.h
usr/include/lirc/media/
usr/include/lirc/media/lirc.h
usr/include/lirc/paths.h
usr/include/lirc/receive.h
usr/include/lirc/release.h
usr/include/lirc/serial.h
usr/include/lirc/transmit.h
usr/include/lirc_client.h
usr/include/lirc_driver.h
usr/include/lirc_private.h
usr/lib/
usr/lib/libirrecord.so
usr/lib/libirrecord.so.0
usr/lib/libirrecord.so.0.0.0
usr/lib/liblirc.so
usr/lib/liblirc.so.0
usr/lib/liblirc.so.0.0.0
usr/lib/liblirc_client.so
usr/lib/liblirc_client.so.0
usr/lib/liblirc_client.so.0.6.0
usr/lib/liblirc_driver.so
usr/lib/liblirc_driver.so.0
usr/lib/liblirc_driver.so.0.3.0
usr/lib/lirc/
usr/lib/lirc/plugins/
usr/lib/lirc/plugins/accent.so
usr/lib/lirc/plugins/alsa_usb.so
usr/lib/lirc/plugins/atilibusb.so
usr/lib/lirc/plugins/atwf83.so
usr/lib/lirc/plugins/audio_alsa.so
usr/lib/lirc/plugins/awlibusb.so
usr/lib/lirc/plugins/bte.so
usr/lib/lirc/plugins/commandir.so
usr/lib/lirc/plugins/creative.so
usr/lib/lirc/plugins/creative_infracd.so
usr/lib/lirc/plugins/default.so
usr/lib/lirc/plugins/devinput.so
usr/lib/lirc/plugins/dfclibusb.so
usr/lib/lirc/plugins/dsp.so
usr/lib/lirc/plugins/ea65.so
usr/lib/lirc/plugins/file.so
usr/lib/lirc/plugins/ftdi.so
usr/lib/lirc/plugins/girs.so
usr/lib/lirc/plugins/hiddev.so
usr/lib/lirc/plugins/i2cuser.so
usr/lib/lirc/plugins/irlink.so
usr/lib/lirc/plugins/irtoy.so
usr/lib/lirc/plugins/livedrive_midi.so
usr/lib/lirc/plugins/livedrive_seq.so
usr/lib/lirc/plugins/logitech.so
usr/lib/lirc/plugins/mouseremote.so
usr/lib/lirc/plugins/mp3anywhere.so
usr/lib/lirc/plugins/mplay.so
usr/lib/lirc/plugins/pcmak.so
usr/lib/lirc/plugins/pinsys.so
usr/lib/lirc/plugins/pixelview.so
usr/lib/lirc/plugins/silitek.so
usr/lib/lirc/plugins/srm7500libusb.so
usr/lib/lirc/plugins/tira.so
usr/lib/lirc/plugins/udp.so
usr/lib/lirc/plugins/uirt2.so
usr/lib/lirc/plugins/uirt2_raw.so
usr/lib/lirc/plugins/usbx.so
usr/lib/lirc/plugins/zotac.so
usr/lib/pkgconfig/
usr/lib/pkgconfig/lirc-driver.pc
usr/lib/pkgconfig/lirc.pc
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/lirc-setup/
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/baseview.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/baseview.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/choosers.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/choosers.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_control.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_control.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_model.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_model.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_view.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/mvc_view.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/util.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc-setup/__pycache__/util.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc-setup/baseview.py
usr/lib/python3.12/site-packages/lirc-setup/choosers.py
usr/lib/python3.12/site-packages/lirc-setup/configs
usr/lib/python3.12/site-packages/lirc-setup/lirc-setup
usr/lib/python3.12/site-packages/lirc-setup/lirc-setup.ui
usr/lib/python3.12/site-packages/lirc-setup/mvc_control.py
usr/lib/python3.12/site-packages/lirc-setup/mvc_model.py
usr/lib/python3.12/site-packages/lirc-setup/mvc_view.py
usr/lib/python3.12/site-packages/lirc-setup/util.py
usr/lib/python3.12/site-packages/lirc/
usr/lib/python3.12/site-packages/lirc/__init__.py
usr/lib/python3.12/site-packages/lirc/__pycache__/
usr/lib/python3.12/site-packages/lirc/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/async_client.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/async_client.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/client.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/client.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/database.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/database.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/paths.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/lirc/__pycache__/paths.cpython-312.pyc
usr/lib/python3.12/site-packages/lirc/_client.so
usr/lib/python3.12/site-packages/lirc/async_client.py
usr/lib/python3.12/site-packages/lirc/client.py
usr/lib/python3.12/site-packages/lirc/database.py
usr/lib/python3.12/site-packages/lirc/paths.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/irexec.service
usr/lib/systemd/system/lircd-setup.service
usr/lib/systemd/system/lircd-uinput.service
usr/lib/systemd/system/lircd.service
usr/lib/systemd/system/lircd.socket
usr/lib/systemd/system/lircmd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/lirc.conf
usr/share/
usr/share/doc/
usr/share/doc/lirc/
usr/share/doc/lirc/VERSION
usr/share/doc/lirc/html/
usr/share/doc/lirc/html/lirc.css
usr/share/doc/lirc/html/table.html
usr/share/doc/lirc/images/
usr/share/doc/lirc/images/atwf83.jpg
usr/share/doc/lirc/images/diode.gif
usr/share/doc/lirc/images/lirc.gif
usr/share/doc/lirc/images/marb18.jpg
usr/share/doc/lirc/irxevent.keys
usr/share/doc/lirc/lirc.org/
usr/share/doc/lirc/lirc.org/PCB.pdf
usr/share/doc/lirc/lirc.org/api-docs
usr/share/doc/lirc/lirc.org/audio-alsa.html
usr/share/doc/lirc/lirc.org/cir.html
usr/share/doc/lirc/lirc.org/cvs.html
usr/share/doc/lirc/lirc.org/developers.html
usr/share/doc/lirc/lirc.org/faq.html
usr/share/doc/lirc/lirc.org/favicon.ico
usr/share/doc/lirc/lirc.org/git.html
usr/share/doc/lirc/lirc.org/html/
usr/share/doc/lirc/lirc.org/html/plugins-index.html
usr/share/doc/lirc/lirc.org/html/table.html
usr/share/doc/lirc/lirc.org/images/
usr/share/doc/lirc/lirc.org/images/IR-mouse1.jpg
usr/share/doc/lirc/lirc.org/images/IR-mouse2.jpg
usr/share/doc/lirc/lirc.org/images/atwf83.jpg
usr/share/doc/lirc/lirc.org/images/audio_transmitter.png
usr/share/doc/lirc/lirc.org/images/audio_wave.png
usr/share/doc/lirc/lirc.org/images/banner_en_1.gif
usr/share/doc/lirc/lirc.org/images/better_transmitter.gif
usr/share/doc/lirc/lirc.org/images/d.jpg
usr/share/doc/lirc/lirc.org/images/diode.gif
usr/share/doc/lirc/lirc.org/images/favicon.ico
usr/share/doc/lirc/lirc.org/images/front.jpg
usr/share/doc/lirc/lirc.org/images/improved_transmitter.gif
usr/share/doc/lirc/lirc.org/images/ir-audio.png
usr/share/doc/lirc/lirc.org/images/irrcv_board.gif
usr/share/doc/lirc/lirc.org/images/lirc.gif
usr/share/doc/lirc/lirc.org/images/lirc1.jpg
usr/share/doc/lirc/lirc.org/images/lirc2.jpg
usr/share/doc/lirc/lirc.org/images/lirc3.jpg
usr/share/doc/lirc/lirc.org/images/lirclogo.gif
usr/share/doc/lirc/lirc.org/images/marb18.jpg
usr/share/doc/lirc/lirc.org/images/schematics.gif
usr/share/doc/lirc/lirc.org/images/screenshot.gif
usr/share/doc/lirc/lirc.org/images/screenshot.jpg
usr/share/doc/lirc/lirc.org/images/screenshot1.gif
usr/share/doc/lirc/lirc.org/images/simple_transmitter.gif
usr/share/doc/lirc/lirc.org/images/top3.jpg
usr/share/doc/lirc/lirc.org/images/tsop173x.png
usr/share/doc/lirc/lirc.org/improved_transmitter.html
usr/share/doc/lirc/lirc.org/index.html
usr/share/doc/lirc/lirc.org/ir-audio.html
usr/share/doc/lirc/lirc.org/irda.html
usr/share/doc/lirc/lirc.org/lirc.css
usr/share/doc/lirc/lirc.org/mirrors.html
usr/share/doc/lirc/lirc.org/parallel.html
usr/share/doc/lirc/lirc.org/pictures.html
usr/share/doc/lirc/lirc.org/pictures1.html
usr/share/doc/lirc/lirc.org/pictures2.html
usr/share/doc/lirc/lirc.org/pinout_TRM1038.txt
usr/share/doc/lirc/lirc.org/receivers.html
usr/share/doc/lirc/lirc.org/robots.txt
usr/share/doc/lirc/lirc.org/software.html
usr/share/doc/lirc/lirc.org/transmitters.html
usr/share/doc/lirc/lirc.org/tv_cards.html
usr/share/doc/lirc/plugindocs/
usr/share/doc/lirc/plugindocs/Makefile
usr/share/doc/lirc/plugindocs/README
usr/share/doc/lirc/plugindocs/driver-toc.xsl
usr/share/doc/lirc/plugindocs/empty_index.tmpl
usr/share/doc/lirc/plugindocs/index.tmpl
usr/share/doc/lirc/plugindocs/lirc.css
usr/share/doc/lirc/plugindocs/make-ext-driver-toc.sh
usr/share/doc/lirc/plugindocs/page.xsl
usr/share/doc/lirc/plugindocs/plugpage.xsl
usr/share/doc/lirc/plugindocs/var
usr/share/lirc/
usr/share/lirc/configs/
usr/share/lirc/configs/accent.conf
usr/share/lirc/configs/act200l.conf
usr/share/lirc/configs/act220l.conf
usr/share/lirc/configs/adaptec.conf
usr/share/lirc/configs/alsa_usb.conf
usr/share/lirc/configs/animax.conf
usr/share/lirc/configs/arduino.conf
usr/share/lirc/configs/asusdh.conf
usr/share/lirc/configs/ati-remote-wonder.conf
usr/share/lirc/configs/atilibusb.conf
usr/share/lirc/configs/atiusb.conf
usr/share/lirc/configs/atwf83.conf
usr/share/lirc/configs/audio.conf
usr/share/lirc/configs/audio_alsa.conf
usr/share/lirc/configs/avermedia.conf
usr/share/lirc/configs/avermedia98.conf
usr/share/lirc/configs/avermedia_vdomate.conf
usr/share/lirc/configs/awlibusb.conf
usr/share/lirc/configs/breakoutbox.conf
usr/share/lirc/configs/bte.conf
usr/share/lirc/configs/chronos.conf
usr/share/lirc/configs/commandir.conf
usr/share/lirc/configs/confs_by_driver.yaml
usr/share/lirc/configs/cph03x_askey.conf
usr/share/lirc/configs/cph03x_dynalink.conf
usr/share/lirc/configs/cph03x_trifan.conf
usr/share/lirc/configs/cph06x_askey.conf
usr/share/lirc/configs/cph06x_phoebe.conf
usr/share/lirc/configs/cph06x_tview.conf
usr/share/lirc/configs/creative.conf
usr/share/lirc/configs/creative_infracd.conf
usr/share/lirc/configs/default.conf
usr/share/lirc/configs/devinput.conf
usr/share/lirc/configs/dfclibusb.conf
usr/share/lirc/configs/drivers.yaml
usr/share/lirc/configs/dsp.conf
usr/share/lirc/configs/dvico.conf
usr/share/lirc/configs/ea65.conf
usr/share/lirc/configs/flyvideo_2.conf
usr/share/lirc/configs/flyvideo_98.conf
usr/share/lirc/configs/flyvideo_98FM.conf
usr/share/lirc/configs/flyvideo_FM_card36.conf
usr/share/lirc/configs/ftdi.conf
usr/share/lirc/configs/gvbctv5pci.conf
usr/share/lirc/configs/hauppauge.conf
usr/share/lirc/configs/hauppauge_dvb.conf
usr/share/lirc/configs/hauppauge_hvr1300.conf
usr/share/lirc/configs/hercules_smarttv_stereo.conf
usr/share/lirc/configs/home_brew_parallel.conf
usr/share/lirc/configs/home_brew_serial.conf
usr/share/lirc/configs/igorplugusb.conf
usr/share/lirc/configs/imon.conf
usr/share/lirc/configs/imon_24g.conf
usr/share/lirc/configs/imon_knob.conf
usr/share/lirc/configs/imon_lcd.conf
usr/share/lirc/configs/imon_pad.conf
usr/share/lirc/configs/imon_rsc.conf
usr/share/lirc/configs/irdeo.conf
usr/share/lirc/configs/irdeo_remote.conf
usr/share/lirc/configs/irlink.conf
usr/share/lirc/configs/irtoy.conf
usr/share/lirc/configs/kernel-drivers.yaml
usr/share/lirc/configs/knc_one.conf
usr/share/lirc/configs/knc_one_typhoon.conf
usr/share/lirc/configs/kworld.conf
usr/share/lirc/configs/leadtek_0007.conf
usr/share/lirc/configs/leadtek_0010.conf
usr/share/lirc/configs/leadtek_pvr2000.conf
usr/share/lirc/configs/livedrive_midi.conf
usr/share/lirc/configs/livedrive_seq.conf
usr/share/lirc/configs/logitech.conf
usr/share/lirc/configs/macmini.conf
usr/share/lirc/configs/mce.conf
usr/share/lirc/configs/mediafocusI.conf
usr/share/lirc/configs/medion.conf
usr/share/lirc/configs/mouseremote.conf
usr/share/lirc/configs/mp3-anywhere-usb.conf
usr/share/lirc/configs/mp3anywhere-serial.conf
usr/share/lirc/configs/mplay_blast.conf
usr/share/lirc/configs/mplay_mini.conf
usr/share/lirc/configs/none.conf
usr/share/lirc/configs/nslu2.conf
usr/share/lirc/configs/nvidia-x10.conf
usr/share/lirc/configs/packard_bell.conf
usr/share/lirc/configs/pcmak.conf
usr/share/lirc/configs/pcmak_usb.conf
usr/share/lirc/configs/pctv.conf
usr/share/lirc/configs/pctv_miro.conf
usr/share/lirc/configs/pctv_pro.conf
usr/share/lirc/configs/pctv_sat.conf
usr/share/lirc/configs/pixelview_bt878.conf
usr/share/lirc/configs/pixelview_mpeg2.conf
usr/share/lirc/configs/pixelview_pak.conf
usr/share/lirc/configs/pixelview_pak_lenco.conf
usr/share/lirc/configs/pixelview_pro_37.conf
usr/share/lirc/configs/pixelview_pv.conf
usr/share/lirc/configs/provideo.conf
usr/share/lirc/configs/realmagic.conf
usr/share/lirc/configs/remotemaster.conf
usr/share/lirc/configs/remotes.list
usr/share/lirc/configs/sa1100_hp.conf
usr/share/lirc/configs/sa1100_sharp.conf
usr/share/lirc/configs/samsung.conf
usr/share/lirc/configs/sasem.conf
usr/share/lirc/configs/sasem_dign.conf
usr/share/lirc/configs/silitek.conf
usr/share/lirc/configs/sir.conf
usr/share/lirc/configs/snapstream.conf
usr/share/lirc/configs/sonyir-osx.conf
usr/share/lirc/configs/sonyir.conf
usr/share/lirc/configs/srm7500libusb.conf
usr/share/lirc/configs/tekram.conf
usr/share/lirc/configs/tekram_bt829.conf
usr/share/lirc/configs/tira.conf
usr/share/lirc/configs/tira_raw.conf
usr/share/lirc/configs/ttusbir.conf
usr/share/lirc/configs/tvbox.conf
usr/share/lirc/configs/udp.conf
usr/share/lirc/configs/uirt2.conf
usr/share/lirc/configs/uirt2_raw.conf
usr/share/lirc/configs/usb_uirt_raw.conf
usr/share/lirc/configs/usbx.conf
usr/share/lirc/configs/x10.conf
usr/share/lirc/configs/zotac.conf
usr/share/lirc/contrib/
usr/share/lirc/contrib/60-lirc.rules
usr/share/lirc/contrib/60-usb-generic-perms.rules
usr/share/lirc/contrib/61-lirc.blacklist-all.conf
usr/share/lirc/contrib/97-lircd-uinput.rules
usr/share/lirc/contrib/98-lirc-stable-link.rules
usr/share/lirc/contrib/99-remote-control-lirc.rules
usr/share/lirc/contrib/irexec.desktop
usr/share/lirc/contrib/irman2lirc
usr/share/lirc/contrib/lircd.conf
usr/share/lirc/contrib/lircmd.conf
usr/share/lirc/contrib/lircrc/
usr/share/lirc/contrib/lircrc/lircmap.xml
usr/share/lirc/contrib/lircrc/lircrc
usr/share/lirc/contrib/lircrc/mplayer.lircrc
usr/share/lirc/contrib/lircrc/mythtv.lircrc
usr/share/lirc/contrib/lircrc/vlc.lircrc
usr/share/lirc/contrib/lircrc/volume_alsa.lircrc
usr/share/lirc/contrib/lircrc/volume_mythtv.lircrc
usr/share/lirc/contrib/lircrc/volume_pulse.lircrc
usr/share/lirc/contrib/release-process.txt
usr/share/lirc/doc/
usr/share/lirc/doc/docpage.xsl
usr/share/lirc/lirc-0.10.2.tar.gz
usr/share/lirc/python-pkg/
usr/share/lirc/python-pkg/README.rst
usr/share/lirc/python-pkg/__pycache__/
usr/share/lirc/python-pkg/__pycache__/config.cpython-312.opt-1.pyc
usr/share/lirc/python-pkg/__pycache__/config.cpython-312.pyc
usr/share/lirc/python-pkg/config.py
usr/share/lirc/python-pkg/lirc/
usr/share/lirc/python-pkg/lirc/__init__.py
usr/share/lirc/python-pkg/lirc/_client.c
usr/share/lirc/python-pkg/lirc/async_client.py
usr/share/lirc/python-pkg/lirc/client.py
usr/share/lirc/python-pkg/lirc/database.py
usr/share/lirc/python-pkg/lirc/doc/
usr/share/lirc/python-pkg/lirc/doc/Doxyfile
usr/share/lirc/python-pkg/lirc/doc/diode-2.gif
usr/share/lirc/python-pkg/lirc/doc/mainpage.dox
usr/share/lirc/python-pkg/lirc/lirctool
usr/share/lirc/python-pkg/lirc/lirctool.py
usr/share/lirc/python-pkg/lirc/paths.py
usr/share/lirc/python-pkg/setup.py
usr/share/lirc/python-pkg/tests/
usr/share/lirc/python-pkg/tests/dummy-server
usr/share/lirc/python-pkg/tests/lircrc.conf
usr/share/lirc/python-pkg/tests/test_client.py
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ircat.1.gz
usr/share/man/man1/irdb-get.1.gz
usr/share/man/man1/irexec.1.gz
usr/share/man/man1/irpipe.1.gz
usr/share/man/man1/irpty.1.gz
usr/share/man/man1/irrecord.1.gz
usr/share/man/man1/irsend.1.gz
usr/share/man/man1/irsimreceive.1.gz
usr/share/man/man1/irsimsend.1.gz
usr/share/man/man1/irtestcase.1.gz
usr/share/man/man1/irtext2udp.1.gz
usr/share/man/man1/irw.1.gz
usr/share/man/man1/irxevent.1.gz
usr/share/man/man1/lirc-config-tool.1.gz
usr/share/man/man1/lirc-lsplugins.1.gz
usr/share/man/man1/lirc-lsremotes.1.gz
usr/share/man/man1/lirc-make-devinput.1.gz
usr/share/man/man1/lirc-postinstall.1.gz
usr/share/man/man1/lirc-setup.1.gz
usr/share/man/man1/mode2.1.gz
usr/share/man/man1/pronto2lirc.1.gz
usr/share/man/man1/xmode2.1.gz
usr/share/man/man5/
usr/share/man/man5/lircd.conf.5.gz
usr/share/man/man5/lircrc.5.gz
usr/share/man/man8/
usr/share/man/man8/lircd-setup.8.gz
usr/share/man/man8/lircd-uinput.8.gz
usr/share/man/man8/lircd.8.gz
usr/share/man/man8/lircmd.8.gz
usr/share/man/man8/lircrcd.8.gz
var/
var/lib/
var/lib/lirc/
var/lib/lirc/images/
var/lib/lirc/images/atwf83.jpg
var/lib/lirc/images/diode.gif
var/lib/lirc/images/lirc.gif
var/lib/lirc/images/marb18.jpg
var/lib/lirc/plugins/
var/lib/lirc/plugins/index.html
var/lib/lirc/plugins/lirc.css
