# extracted/adguardhome-1:0.107.55-1/files
%FILES%
etc/
etc/adguardhome.yaml
usr/
usr/bin/
usr/bin/adguardhome
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/adguardhome.service
