# extracted/prometheus-blackbox-exporter-0.25.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-blackbox-exporter
usr/
usr/bin/
usr/bin/prometheus-blackbox-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-blackbox-exporter.service
