# extracted/hostapd-2.11-4/files
%FILES%
etc/
etc/hostapd/
etc/hostapd/hostapd.accept
etc/hostapd/hostapd.conf
etc/hostapd/hostapd.deny
etc/hostapd/hostapd.eap_user
etc/hostapd/hostapd.radius_clients
etc/hostapd/hostapd.vlan
etc/hostapd/hostapd.wpa_psk
usr/
usr/bin/
usr/bin/hostapd
usr/bin/hostapd_cli
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hostapd.service
usr/lib/systemd/system/hostapd@.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/hostapd.conf
usr/share/
usr/share/doc/
usr/share/doc/hostapd/
usr/share/doc/hostapd/ChangeLog
usr/share/doc/hostapd/README
usr/share/doc/hostapd/README-MULTI-AP
usr/share/doc/hostapd/README-WPS
usr/share/doc/hostapd/hlr_auc_gw.milenage_db
usr/share/doc/hostapd/hlr_auc_gw.txt
usr/share/doc/hostapd/hostapd.sim_db
usr/share/doc/hostapd/wired.conf
usr/share/factory/
usr/share/factory/etc/
usr/share/factory/etc/hostapd/
usr/share/factory/etc/hostapd/hostapd.accept
usr/share/factory/etc/hostapd/hostapd.conf
usr/share/factory/etc/hostapd/hostapd.deny
usr/share/factory/etc/hostapd/hostapd.eap_user
usr/share/factory/etc/hostapd/hostapd.radius_clients
usr/share/factory/etc/hostapd/hostapd.vlan
usr/share/factory/etc/hostapd/hostapd.wpa_psk
usr/share/licenses/
usr/share/licenses/hostapd/
usr/share/licenses/hostapd/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hostapd_cli.1.gz
usr/share/man/man8/
usr/share/man/man8/hostapd.8.gz
