# extracted/lighttpd-1.4.76-1/files
%FILES%
etc/
etc/lighttpd/
etc/lighttpd/lighttpd.conf
etc/logrotate.d/
etc/logrotate.d/lighttpd
usr/
usr/bin/
usr/bin/lighttpd
usr/bin/lighttpd-angel
usr/lib/
usr/lib/lighttpd/
usr/lib/lighttpd/mod_accesslog.so
usr/lib/lighttpd/mod_ajp13.so
usr/lib/lighttpd/mod_auth.so
usr/lib/lighttpd/mod_authn_dbi.so
usr/lib/lighttpd/mod_authn_file.so
usr/lib/lighttpd/mod_authn_gssapi.so
usr/lib/lighttpd/mod_authn_ldap.so
usr/lib/lighttpd/mod_cgi.so
usr/lib/lighttpd/mod_deflate.so
usr/lib/lighttpd/mod_dirlisting.so
usr/lib/lighttpd/mod_extforward.so
usr/lib/lighttpd/mod_h2.so
usr/lib/lighttpd/mod_magnet.so
usr/lib/lighttpd/mod_maxminddb.so
usr/lib/lighttpd/mod_openssl.so
usr/lib/lighttpd/mod_proxy.so
usr/lib/lighttpd/mod_rrdtool.so
usr/lib/lighttpd/mod_sockproxy.so
usr/lib/lighttpd/mod_ssi.so
usr/lib/lighttpd/mod_status.so
usr/lib/lighttpd/mod_userdir.so
usr/lib/lighttpd/mod_vhostdb.so
usr/lib/lighttpd/mod_vhostdb_dbi.so
usr/lib/lighttpd/mod_vhostdb_ldap.so
usr/lib/lighttpd/mod_vhostdb_mysql.so
usr/lib/lighttpd/mod_webdav.so
usr/lib/lighttpd/mod_wstunnel.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lighttpd.service
usr/lib/systemd/system/lighttpd.socket
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/lighttpd.conf
usr/share/
usr/share/doc/
usr/share/doc/lighttpd/
usr/share/doc/lighttpd/config/
usr/share/doc/lighttpd/config/conf.d/
usr/share/doc/lighttpd/config/conf.d/access_log.conf
usr/share/doc/lighttpd/config/conf.d/auth.conf
usr/share/doc/lighttpd/config/conf.d/cgi.conf
usr/share/doc/lighttpd/config/conf.d/debug.conf
usr/share/doc/lighttpd/config/conf.d/deflate.conf
usr/share/doc/lighttpd/config/conf.d/dirlisting.conf
usr/share/doc/lighttpd/config/conf.d/evhost.conf
usr/share/doc/lighttpd/config/conf.d/expire.conf
usr/share/doc/lighttpd/config/conf.d/fastcgi.conf
usr/share/doc/lighttpd/config/conf.d/magnet.conf
usr/share/doc/lighttpd/config/conf.d/mime.conf
usr/share/doc/lighttpd/config/conf.d/mod.template
usr/share/doc/lighttpd/config/conf.d/proxy.conf
usr/share/doc/lighttpd/config/conf.d/rrdtool.conf
usr/share/doc/lighttpd/config/conf.d/scgi.conf
usr/share/doc/lighttpd/config/conf.d/simple_vhost.conf
usr/share/doc/lighttpd/config/conf.d/ssi.conf
usr/share/doc/lighttpd/config/conf.d/status.conf
usr/share/doc/lighttpd/config/conf.d/userdir.conf
usr/share/doc/lighttpd/config/conf.d/webdav.conf
usr/share/doc/lighttpd/config/lighttpd.conf
usr/share/doc/lighttpd/config/modules.conf
usr/share/doc/lighttpd/config/vhosts.d/
usr/share/doc/lighttpd/config/vhosts.d/vhosts.template
usr/share/licenses/
usr/share/licenses/lighttpd/
usr/share/licenses/lighttpd/COPYING
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/lighttpd-angel.8.gz
usr/share/man/man8/lighttpd.8.gz
var/
var/cache/
var/cache/lighttpd/
var/log/
var/log/lighttpd/
