# extracted/zabbix-server-7.2.0-1/files
%FILES%
etc/
etc/zabbix/
etc/zabbix/zabbix_server.conf
usr/
usr/bin/
usr/bin/zabbix_get
usr/bin/zabbix_server_mysql
usr/bin/zabbix_server_postgresql
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/zabbix-server-mysql.service
usr/lib/systemd/system/zabbix-server-pgsql.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/zabbix-server.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/zabbix-server.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/zabbix_get.1.gz
usr/share/man/man8/
usr/share/man/man8/zabbix_server.8.gz
usr/share/zabbix-server/
usr/share/zabbix-server/mysql/
usr/share/zabbix-server/mysql/data.sql
usr/share/zabbix-server/mysql/images.sql
usr/share/zabbix-server/mysql/schema.sql
usr/share/zabbix-server/postgresql/
usr/share/zabbix-server/postgresql/data.sql
usr/share/zabbix-server/postgresql/images.sql
usr/share/zabbix-server/postgresql/schema.sql
