# extracted/openfortivpn-1.22.1-1/files
%FILES%
etc/
etc/openfortivpn/
etc/openfortivpn/config
usr/
usr/bin/
usr/bin/openfortivpn
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openfortivpn@.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/openfortivpn.1.gz
usr/share/openfortivpn/
usr/share/openfortivpn/config.template
