# extracted/apcupsd-3.14.14-8/files
%FILES%
etc/
etc/apcupsd/
etc/apcupsd/apccontrol
etc/apcupsd/apcupsd.conf
etc/apcupsd/apcupsd.css
etc/apcupsd/changeme
etc/apcupsd/commfailure
etc/apcupsd/commok
etc/apcupsd/hosts.conf
etc/apcupsd/multimon.conf
etc/apcupsd/offbattery
etc/apcupsd/onbattery
usr/
usr/bin/
usr/bin/apcaccess
usr/bin/apctest
usr/bin/apcupsd
usr/bin/smtp-apcupsd
usr/lib/
usr/lib/apcupsd/
usr/lib/apcupsd/cgi-bin/
usr/lib/apcupsd/cgi-bin/multimon.cgi
usr/lib/apcupsd/cgi-bin/upsfstats.cgi
usr/lib/apcupsd/cgi-bin/upsimage.cgi
usr/lib/apcupsd/cgi-bin/upsstats.cgi
usr/lib/systemd/
usr/lib/systemd/system-shutdown/
usr/lib/systemd/system-shutdown/apcupsd_shutdown
usr/lib/systemd/system/
usr/lib/systemd/system/apcupsd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/apcupsd.conf
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/apcupsd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/apcaccess.8.gz
usr/share/man/man8/apccontrol.8.gz
usr/share/man/man8/apctest.8.gz
usr/share/man/man8/apcupsd.8.gz
