# extracted/zabbix-agent-7.2.0-1/files
%FILES%
etc/
etc/zabbix/
etc/zabbix/zabbix_agentd.conf
usr/
usr/bin/
usr/bin/zabbix_agentd
usr/bin/zabbix_sender
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/zabbix-agent.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/zabbix-agent.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/zabbix-agent.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/zabbix_sender.1.gz
usr/share/man/man8/
usr/share/man/man8/zabbix_agentd.8.gz
usr/share/zabbix-agent/
usr/share/zabbix-agent/userparameter_examples.conf
usr/share/zabbix-agent/userparameter_mysql.conf
