# extracted/gpsd-3.25-5/files
%FILES%
etc/
etc/default/
etc/default/gpsd
usr/
usr/bin/
usr/bin/cgps
usr/bin/gegps
usr/bin/gps2udp
usr/bin/gpscat
usr/bin/gpscsv
usr/bin/gpsctl
usr/bin/gpsd
usr/bin/gpsdctl
usr/bin/gpsdebuginfo
usr/bin/gpsdecode
usr/bin/gpsfake
usr/bin/gpsinit
usr/bin/gpsmon
usr/bin/gpspipe
usr/bin/gpsplot
usr/bin/gpsprof
usr/bin/gpsrinex
usr/bin/gpssnmp
usr/bin/gpssubframe
usr/bin/gpxlogger
usr/bin/lcdgps
usr/bin/ntpshmmon
usr/bin/ppscheck
usr/bin/ubxtool
usr/bin/xgps
usr/bin/xgpsspeed
usr/bin/zerk
usr/include/
usr/include/gps.h
usr/include/libgpsmm.h
usr/lib/
usr/lib/libQgpsmm.prl
usr/lib/libQgpsmm.so
usr/lib/libQgpsmm.so.30
usr/lib/libQgpsmm.so.30.0.0
usr/lib/libgps.so
usr/lib/libgps.so.30
usr/lib/libgps.so.30.0.0
usr/lib/libgpsdpacket.so
usr/lib/libgpsdpacket.so.30
usr/lib/libgpsdpacket.so.30.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/Qgpsmm.pc
usr/lib/pkgconfig/libgps.pc
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/gps-3.25.egg-info
usr/lib/python3.13/site-packages/gps/
usr/lib/python3.13/site-packages/gps/__init__.py
usr/lib/python3.13/site-packages/gps/aiogps.py
usr/lib/python3.13/site-packages/gps/client.py
usr/lib/python3.13/site-packages/gps/clienthelpers.py
usr/lib/python3.13/site-packages/gps/fake.py
usr/lib/python3.13/site-packages/gps/gps.py
usr/lib/python3.13/site-packages/gps/misc.py
usr/lib/python3.13/site-packages/gps/packet.py
usr/lib/python3.13/site-packages/gps/ubx.py
usr/lib/python3.13/site-packages/gps/watch_options.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gpsd.service
usr/lib/systemd/system/gpsd.socket
usr/lib/systemd/system/gpsdctl@.service
usr/lib/udev/
usr/lib/udev/gpsd.hotplug
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/25-gpsd-usb.rules
usr/share/
usr/share/applications/
usr/share/applications/xgps.desktop
usr/share/applications/xgpsspeed.desktop
usr/share/gpsd/
usr/share/gpsd/doc/
usr/share/gpsd/doc/AUTHORS
usr/share/gpsd/doc/COPYING
usr/share/gpsd/doc/NEWS
usr/share/gpsd/doc/README.adoc
usr/share/gpsd/doc/SUPPORT.adoc
usr/share/gpsd/doc/build.adoc
usr/share/gpsd/doc/example1.c.txt
usr/share/gpsd/doc/example2.py.txt
usr/share/gpsd/icons/
usr/share/gpsd/icons/gpsd-logo.png
usr/share/licenses/
usr/share/licenses/gpsd/
usr/share/licenses/gpsd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/cgps.1.gz
usr/share/man/man1/gegps.1.gz
usr/share/man/man1/gps.1.gz
usr/share/man/man1/gps2udp.1.gz
usr/share/man/man1/gpscat.1.gz
usr/share/man/man1/gpscsv.1.gz
usr/share/man/man1/gpsctl.1.gz
usr/share/man/man1/gpsdebuginfo.1.gz
usr/share/man/man1/gpsdecode.1.gz
usr/share/man/man1/gpsfake.1.gz
usr/share/man/man1/gpsmon.1.gz
usr/share/man/man1/gpspipe.1.gz
usr/share/man/man1/gpsplot.1.gz
usr/share/man/man1/gpsprof.1.gz
usr/share/man/man1/gpsrinex.1.gz
usr/share/man/man1/gpssnmp.1.gz
usr/share/man/man1/gpssubframe.1.gz
usr/share/man/man1/gpxlogger.1.gz
usr/share/man/man1/lcdgps.1.gz
usr/share/man/man1/ntpshmmon.1.gz
usr/share/man/man1/ubxtool.1.gz
usr/share/man/man1/xgps.1.gz
usr/share/man/man1/xgpsspeed.1.gz
usr/share/man/man1/zerk.1.gz
usr/share/man/man3/
usr/share/man/man3/libQgpsmm.3.gz
usr/share/man/man3/libgps.3.gz
usr/share/man/man3/libgpsmm.3.gz
usr/share/man/man5/
usr/share/man/man5/gpsd_json.5.gz
usr/share/man/man8/
usr/share/man/man8/gpsd.8.gz
usr/share/man/man8/gpsdctl.8.gz
usr/share/man/man8/gpsinit.8.gz
usr/share/man/man8/ppscheck.8.gz
usr/share/snmp/
usr/share/snmp/mibs/
usr/share/snmp/mibs/gpsd/
usr/share/snmp/mibs/gpsd/GPSD-MIB
