# extracted/python-geoipsets-2.4-2/files
%FILES%
usr/
usr/bin/
usr/bin/geoipsets
usr/lib/
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/PKG-INFO
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/SOURCES.txt
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/dependency_links.txt
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/entry_points.txt
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/requires.txt
usr/lib/python3.13/site-packages/geoipsets-2.4.0-py3.13.egg-info/top_level.txt
usr/lib/python3.13/site-packages/geoipsets/
usr/lib/python3.13/site-packages/geoipsets/VERSION
usr/lib/python3.13/site-packages/geoipsets/__init__.py
usr/lib/python3.13/site-packages/geoipsets/__main__.py
usr/lib/python3.13/site-packages/geoipsets/__pycache__/
usr/lib/python3.13/site-packages/geoipsets/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/__main__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/__main__.cpython-313.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/dbip.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/dbip.cpython-313.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/maxmind.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/maxmind.cpython-313.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/utils.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/geoipsets/__pycache__/utils.cpython-313.pyc
usr/lib/python3.13/site-packages/geoipsets/dbip.py
usr/lib/python3.13/site-packages/geoipsets/maxmind.py
usr/lib/python3.13/site-packages/geoipsets/utils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/update-geoipsets.service
usr/lib/systemd/system/update-geoipsets.timer
usr/share/
usr/share/doc/
usr/share/doc/python-geoipsets/
usr/share/doc/python-geoipsets/README.md
usr/share/doc/python-geoipsets/README_python.md
usr/share/licenses/
usr/share/licenses/python-geoipsets/
usr/share/licenses/python-geoipsets/LICENSE
