# extracted/mirro-rs-0.2.3-1/files
%FILES%
usr/
usr/bin/
usr/bin/mirro-rs
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mirro-rs.service
usr/lib/systemd/system/mirro-rs.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/mirro-rs.bash
usr/share/doc/
usr/share/doc/mirro-rs/
usr/share/doc/mirro-rs/README.md
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/mirro-rs.fish
usr/share/licenses/
usr/share/licenses/mirro-rs/
usr/share/licenses/mirro-rs/LICENSE-MIT
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/_mirro-rs.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_mirro-rs
