# extracted/deepin-app-services-1.0.26-2/files
%FILES%
usr/
usr/bin/
usr/bin/dde-dconfig
usr/bin/dde-dconfig-daemon
usr/bin/dde-dconfig-editor
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dde-dconfig-daemon.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/dde-dconfig-daemon.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/dde-dconfig
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.desktopspec.ConfigManager.Manager.xml
usr/share/dbus-1/interfaces/org.desktopspec.ConfigManager.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.desktopspec.ConfigManager.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.desktopspec.ConfigManager.conf
usr/share/dde-dconfig-editor/
usr/share/dde-dconfig-editor/translations/
usr/share/dde-dconfig-editor/translations/dde-dconfig-editor_zh_CN.qm
usr/share/dde-dconfig/
usr/share/dde-dconfig/translations/
usr/share/dde-dconfig/translations/dde-dconfig_zh_CN.qm
usr/share/dsg/
usr/share/dsg/configs/
usr/share/dsg/configs/dconfig-example/
usr/share/dsg/configs/dconfig-example/a/
usr/share/dsg/configs/dconfig-example/a/example.json
usr/share/dsg/configs/dconfig-example/example.json
usr/share/dsg/configs/example.json
usr/share/dsg/configs/overrides/
usr/share/dsg/configs/overrides/dconfig-example/
usr/share/dsg/configs/overrides/dconfig-example/example/
usr/share/dsg/configs/overrides/dconfig-example/example/a/
usr/share/dsg/configs/overrides/dconfig-example/example/a/dconf-example.override.a.json
usr/share/dsg/configs/overrides/dconfig-example/example/dconf-example.override.a.json
usr/share/dsg/configs/overrides/dconfig-example/example/dconf-example.override.json
usr/share/qt5/
usr/share/qt5/doc/
usr/share/qt5/doc/dde-dconfig-doc.qch
usr/share/zsh/
usr/share/zsh/vendor-completions/
usr/share/zsh/vendor-completions/_dde-dconfig
