# extracted/fwknop-2.6.11-1/files
%FILES%
etc/
etc/fwknop/
etc/fwknop/access.conf
etc/fwknop/fwknopd.conf
usr/
usr/bin/
usr/bin/fwknop
usr/bin/fwknopd
usr/include/
usr/include/fko.h
usr/lib/
usr/lib/libfko.so
usr/lib/libfko.so.3
usr/lib/libfko.so.3.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/fwknopd.service
usr/share/
usr/share/info/
usr/share/info/libfko.info.gz
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/fwknop.8.gz
usr/share/man/man8/fwknopd.8.gz
