# extracted/gitlab-pages-17.6.2-1/files
%FILES%
etc/
etc/gitlab-pages/
etc/gitlab-pages/config.cfg
etc/gitlab-pages/service.env
usr/
usr/bin/
usr/bin/gitlab-pages
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gitlab-pages.service
