# extracted/grub-btrfs-4.13-2/files
%FILES%
etc/
etc/default/
etc/default/grub-btrfs/
etc/default/grub-btrfs/config
etc/grub.d/
etc/grub.d/41_snapshots-btrfs
usr/
usr/bin/
usr/bin/grub-btrfsd
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/grub-btrfs-overlayfs
usr/lib/initcpio/install/
usr/lib/initcpio/install/grub-btrfs-overlayfs
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/grub-btrfsd.service
usr/share/
usr/share/doc/
usr/share/doc/grub-btrfs/
usr/share/doc/grub-btrfs/README.md
usr/share/doc/grub-btrfs/initramfs-overlayfs.md
usr/share/licenses/
usr/share/licenses/grub-btrfs/
usr/share/licenses/grub-btrfs/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/grub-btrfs.8.bz2
usr/share/man/man8/grub-btrfsd.8.bz2
