# extracted/postgrey-1.37-6/files
%FILES%
etc/
etc/postfix/
etc/postfix/postgrey_whitelist_clients
etc/postfix/postgrey_whitelist_recipients
usr/
usr/bin/
usr/bin/postgrey
usr/bin/postgreyreport
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/postgrey.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/postgrey.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/postgrey.conf
