# extracted/ossp-1.3.2+11+g9e94d67-4/files
%FILES%
usr/
usr/bin/
usr/bin/ossp-alsap
usr/bin/ossp-padsp
usr/bin/osspd
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/osspd.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/osspd.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/98-osscuse.rules
usr/share/
usr/share/doc/
usr/share/doc/ossp/
usr/share/doc/ossp/README
