# extracted/zeroc-ice-3.7.10-1/files
%FILES%
etc/
etc/glacier2router.conf
etc/icegridnode.conf
etc/icegridregistry.conf
usr/
usr/bin/
usr/bin/glacier2router
usr/bin/icebox
usr/bin/icebox++11
usr/bin/iceboxadmin
usr/bin/icebridge
usr/bin/icegridadmin
usr/bin/icegriddb
usr/bin/icegridnode
usr/bin/icegridregistry
usr/bin/icepatch2calc
usr/bin/icepatch2client
usr/bin/icepatch2server
usr/bin/icestormadmin
usr/bin/icestormdb
usr/bin/slice2confluence
usr/bin/slice2cpp
usr/bin/slice2cs
usr/bin/slice2html
usr/bin/slice2java
usr/bin/slice2js
usr/bin/slice2matlab
usr/bin/slice2objc
usr/bin/slice2php
usr/bin/slice2rb
usr/bin/slice2swift
usr/include/
usr/include/Glacier2/
usr/include/Glacier2/Application.h
usr/include/Glacier2/Config.h
usr/include/Glacier2/Glacier2.h
usr/include/Glacier2/Metrics.h
usr/include/Glacier2/NullPermissionsVerifier.h
usr/include/Glacier2/PermissionsVerifier.h
usr/include/Glacier2/PermissionsVerifierF.h
usr/include/Glacier2/Router.h
usr/include/Glacier2/RouterF.h
usr/include/Glacier2/SSLInfo.h
usr/include/Glacier2/Session.h
usr/include/Glacier2/SessionHelper.h
usr/include/Ice/
usr/include/Ice/Application.h
usr/include/Ice/AsyncResult.h
usr/include/Ice/AsyncResultF.h
usr/include/Ice/BatchRequestInterceptor.h
usr/include/Ice/BatchRequestQueueF.h
usr/include/Ice/Buffer.h
usr/include/Ice/BuiltinSequences.h
usr/include/Ice/Communicator.h
usr/include/Ice/CommunicatorAsync.h
usr/include/Ice/CommunicatorF.h
usr/include/Ice/Comparable.h
usr/include/Ice/Config.h
usr/include/Ice/Connection.h
usr/include/Ice/ConnectionAsync.h
usr/include/Ice/ConnectionF.h
usr/include/Ice/ConnectionIF.h
usr/include/Ice/ConsoleUtil.h
usr/include/Ice/Current.h
usr/include/Ice/DefaultValueFactory.h
usr/include/Ice/DispatchInterceptor.h
usr/include/Ice/Dispatcher.h
usr/include/Ice/DynamicLibrary.h
usr/include/Ice/DynamicLibraryF.h
usr/include/Ice/Endpoint.h
usr/include/Ice/EndpointF.h
usr/include/Ice/EndpointTypes.h
usr/include/Ice/Exception.h
usr/include/Ice/ExceptionHelpers.h
usr/include/Ice/FacetMap.h
usr/include/Ice/FactoryTable.h
usr/include/Ice/FactoryTableInit.h
usr/include/Ice/Format.h
usr/include/Ice/Functional.h
usr/include/Ice/GCObject.h
usr/include/Ice/Handle.h
usr/include/Ice/Ice.h
usr/include/Ice/IconvStringConverter.h
usr/include/Ice/Identity.h
usr/include/Ice/ImplicitContext.h
usr/include/Ice/ImplicitContextF.h
usr/include/Ice/Incoming.h
usr/include/Ice/IncomingAsync.h
usr/include/Ice/IncomingAsyncF.h
usr/include/Ice/Initialize.h
usr/include/Ice/InputStream.h
usr/include/Ice/InstanceF.h
usr/include/Ice/Instrumentation.h
usr/include/Ice/InstrumentationF.h
usr/include/Ice/InterfaceByValue.h
usr/include/Ice/LocalException.h
usr/include/Ice/LocalObject.h
usr/include/Ice/LocalObjectF.h
usr/include/Ice/Locator.h
usr/include/Ice/LocatorF.h
usr/include/Ice/Logger.h
usr/include/Ice/LoggerF.h
usr/include/Ice/LoggerUtil.h
usr/include/Ice/Metrics.h
usr/include/Ice/MetricsAdminI.h
usr/include/Ice/MetricsFunctional.h
usr/include/Ice/MetricsObserverI.h
usr/include/Ice/NativePropertiesAdmin.h
usr/include/Ice/Object.h
usr/include/Ice/ObjectAdapter.h
usr/include/Ice/ObjectAdapterF.h
usr/include/Ice/ObjectF.h
usr/include/Ice/ObjectFactory.h
usr/include/Ice/ObserverHelper.h
usr/include/Ice/Optional.h
usr/include/Ice/OutgoingAsync.h
usr/include/Ice/OutgoingAsyncF.h
usr/include/Ice/OutputStream.h
usr/include/Ice/Plugin.h
usr/include/Ice/PluginF.h
usr/include/Ice/Process.h
usr/include/Ice/ProcessF.h
usr/include/Ice/Properties.h
usr/include/Ice/PropertiesAdmin.h
usr/include/Ice/PropertiesF.h
usr/include/Ice/Protocol.h
usr/include/Ice/Proxy.h
usr/include/Ice/ProxyF.h
usr/include/Ice/ProxyFactoryF.h
usr/include/Ice/ProxyHandle.h
usr/include/Ice/ReferenceF.h
usr/include/Ice/RegisterPlugins.h
usr/include/Ice/RemoteLogger.h
usr/include/Ice/RequestHandlerF.h
usr/include/Ice/ResponseHandlerF.h
usr/include/Ice/Router.h
usr/include/Ice/RouterF.h
usr/include/Ice/SHA1.h
usr/include/Ice/ServantLocator.h
usr/include/Ice/ServantLocatorF.h
usr/include/Ice/ServantManagerF.h
usr/include/Ice/Service.h
usr/include/Ice/SliceChecksumDict.h
usr/include/Ice/SliceChecksums.h
usr/include/Ice/SlicedData.h
usr/include/Ice/SlicedDataF.h
usr/include/Ice/StreamHelpers.h
usr/include/Ice/StringConverter.h
usr/include/Ice/ThreadPoolF.h
usr/include/Ice/UUID.h
usr/include/Ice/UniquePtr.h
usr/include/Ice/UniqueRef.h
usr/include/Ice/UserExceptionFactory.h
usr/include/Ice/Value.h
usr/include/Ice/ValueF.h
usr/include/Ice/ValueFactory.h
usr/include/Ice/Version.h
usr/include/IceBox/
usr/include/IceBox/Config.h
usr/include/IceBox/IceBox.h
usr/include/IceGrid/
usr/include/IceGrid/Admin.h
usr/include/IceGrid/Config.h
usr/include/IceGrid/Descriptor.h
usr/include/IceGrid/Exception.h
usr/include/IceGrid/FileParser.h
usr/include/IceGrid/IceGrid.h
usr/include/IceGrid/PluginFacade.h
usr/include/IceGrid/Registry.h
usr/include/IceGrid/ReplicaGroupFilter.h
usr/include/IceGrid/Session.h
usr/include/IceGrid/UserAccountMapper.h
usr/include/IcePatch2/
usr/include/IcePatch2/ClientUtil.h
usr/include/IcePatch2/Config.h
usr/include/IcePatch2/FileInfo.h
usr/include/IcePatch2/FileServer.h
usr/include/IcePatch2/IcePatch2.h
usr/include/IceSSL/
usr/include/IceSSL/Config.h
usr/include/IceSSL/ConnectionInfo.h
usr/include/IceSSL/ConnectionInfoF.h
usr/include/IceSSL/EndpointInfo.h
usr/include/IceSSL/IceSSL.h
usr/include/IceSSL/OpenSSL.h
usr/include/IceSSL/Plugin.h
usr/include/IceSSL/SChannel.h
usr/include/IceSSL/SecureTransport.h
usr/include/IceStorm/
usr/include/IceStorm/Config.h
usr/include/IceStorm/IceStorm.h
usr/include/IceStorm/Metrics.h
usr/include/IceUtil/
usr/include/IceUtil/Atomic.h
usr/include/IceUtil/Cond.h
usr/include/IceUtil/Config.h
usr/include/IceUtil/ConsoleUtil.h
usr/include/IceUtil/CountDownLatch.h
usr/include/IceUtil/CtrlCHandler.h
usr/include/IceUtil/DisableWarnings.h
usr/include/IceUtil/Exception.h
usr/include/IceUtil/FileUtil.h
usr/include/IceUtil/Functional.h
usr/include/IceUtil/Handle.h
usr/include/IceUtil/IceUtil.h
usr/include/IceUtil/InputUtil.h
usr/include/IceUtil/Iterator.h
usr/include/IceUtil/Lock.h
usr/include/IceUtil/Monitor.h
usr/include/IceUtil/Mutex.h
usr/include/IceUtil/MutexProtocol.h
usr/include/IceUtil/MutexPtrLock.h
usr/include/IceUtil/MutexPtrTryLock.h
usr/include/IceUtil/Optional.h
usr/include/IceUtil/Options.h
usr/include/IceUtil/OutputUtil.h
usr/include/IceUtil/PopDisableWarnings.h
usr/include/IceUtil/PushDisableWarnings.h
usr/include/IceUtil/Random.h
usr/include/IceUtil/RecMutex.h
usr/include/IceUtil/ResourceConfig.h
usr/include/IceUtil/ScannerConfig.h
usr/include/IceUtil/ScopedArray.h
usr/include/IceUtil/Shared.h
usr/include/IceUtil/StopWatch.h
usr/include/IceUtil/StringConverter.h
usr/include/IceUtil/StringUtil.h
usr/include/IceUtil/Thread.h
usr/include/IceUtil/ThreadException.h
usr/include/IceUtil/Time.h
usr/include/IceUtil/Timer.h
usr/include/IceUtil/UUID.h
usr/include/IceUtil/UndefSysMacros.h
usr/lib/
usr/lib/libGlacier2++11.so
usr/lib/libGlacier2++11.so.3.7.10
usr/lib/libGlacier2++11.so.37
usr/lib/libGlacier2.so
usr/lib/libGlacier2.so.3.7.10
usr/lib/libGlacier2.so.37
usr/lib/libGlacier2CryptPermissionsVerifier.so.3.7.10
usr/lib/libGlacier2CryptPermissionsVerifier.so.37
usr/lib/libIce++11.so
usr/lib/libIce++11.so.3.7.10
usr/lib/libIce++11.so.37
usr/lib/libIce.so
usr/lib/libIce.so.3.7.10
usr/lib/libIce.so.37
usr/lib/libIceBox++11.so
usr/lib/libIceBox++11.so.3.7.10
usr/lib/libIceBox++11.so.37
usr/lib/libIceBox.so
usr/lib/libIceBox.so.3.7.10
usr/lib/libIceBox.so.37
usr/lib/libIceDB.so.3.7.10
usr/lib/libIceDB.so.37
usr/lib/libIceDiscovery++11.so
usr/lib/libIceDiscovery++11.so.3.7.10
usr/lib/libIceDiscovery++11.so.37
usr/lib/libIceDiscovery.so
usr/lib/libIceDiscovery.so.3.7.10
usr/lib/libIceDiscovery.so.37
usr/lib/libIceGrid++11.so
usr/lib/libIceGrid++11.so.3.7.10
usr/lib/libIceGrid++11.so.37
usr/lib/libIceGrid.so
usr/lib/libIceGrid.so.3.7.10
usr/lib/libIceGrid.so.37
usr/lib/libIceLocatorDiscovery++11.so
usr/lib/libIceLocatorDiscovery++11.so.3.7.10
usr/lib/libIceLocatorDiscovery++11.so.37
usr/lib/libIceLocatorDiscovery.so
usr/lib/libIceLocatorDiscovery.so.3.7.10
usr/lib/libIceLocatorDiscovery.so.37
usr/lib/libIcePatch2.so
usr/lib/libIcePatch2.so.3.7.10
usr/lib/libIcePatch2.so.37
usr/lib/libIceSSL++11.so
usr/lib/libIceSSL++11.so.3.7.10
usr/lib/libIceSSL++11.so.37
usr/lib/libIceSSL.so
usr/lib/libIceSSL.so.3.7.10
usr/lib/libIceSSL.so.37
usr/lib/libIceStorm++11.so
usr/lib/libIceStorm++11.so.3.7.10
usr/lib/libIceStorm++11.so.37
usr/lib/libIceStorm.so
usr/lib/libIceStorm.so.3.7.10
usr/lib/libIceStorm.so.37
usr/lib/libIceStormService.so.3.7.10
usr/lib/libIceStormService.so.37
usr/lib/libIceXML.so.3.7.10
usr/lib/libIceXML.so.37
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/glacier2router.service
usr/lib/systemd/system/icegridnode.service
usr/lib/systemd/system/icegridregistry.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/zeroc-ice.conf
usr/share/
usr/share/ice/
usr/share/ice/ICE_LICENSE
usr/share/ice/LICENSE
usr/share/ice/slice/
usr/share/ice/slice/Glacier2/
usr/share/ice/slice/Glacier2/Metrics.ice
usr/share/ice/slice/Glacier2/PermissionsVerifier.ice
usr/share/ice/slice/Glacier2/PermissionsVerifierF.ice
usr/share/ice/slice/Glacier2/Router.ice
usr/share/ice/slice/Glacier2/RouterF.ice
usr/share/ice/slice/Glacier2/SSLInfo.ice
usr/share/ice/slice/Glacier2/Session.ice
usr/share/ice/slice/Ice/
usr/share/ice/slice/Ice/BuiltinSequences.ice
usr/share/ice/slice/Ice/Communicator.ice
usr/share/ice/slice/Ice/CommunicatorF.ice
usr/share/ice/slice/Ice/Connection.ice
usr/share/ice/slice/Ice/ConnectionF.ice
usr/share/ice/slice/Ice/Current.ice
usr/share/ice/slice/Ice/Endpoint.ice
usr/share/ice/slice/Ice/EndpointF.ice
usr/share/ice/slice/Ice/EndpointTypes.ice
usr/share/ice/slice/Ice/FacetMap.ice
usr/share/ice/slice/Ice/Identity.ice
usr/share/ice/slice/Ice/ImplicitContext.ice
usr/share/ice/slice/Ice/ImplicitContextF.ice
usr/share/ice/slice/Ice/Instrumentation.ice
usr/share/ice/slice/Ice/InstrumentationF.ice
usr/share/ice/slice/Ice/LocalException.ice
usr/share/ice/slice/Ice/Locator.ice
usr/share/ice/slice/Ice/LocatorF.ice
usr/share/ice/slice/Ice/Logger.ice
usr/share/ice/slice/Ice/LoggerF.ice
usr/share/ice/slice/Ice/Metrics.ice
usr/share/ice/slice/Ice/ObjectAdapter.ice
usr/share/ice/slice/Ice/ObjectAdapterF.ice
usr/share/ice/slice/Ice/ObjectFactory.ice
usr/share/ice/slice/Ice/Plugin.ice
usr/share/ice/slice/Ice/PluginF.ice
usr/share/ice/slice/Ice/Process.ice
usr/share/ice/slice/Ice/ProcessF.ice
usr/share/ice/slice/Ice/Properties.ice
usr/share/ice/slice/Ice/PropertiesAdmin.ice
usr/share/ice/slice/Ice/PropertiesF.ice
usr/share/ice/slice/Ice/RemoteLogger.ice
usr/share/ice/slice/Ice/Router.ice
usr/share/ice/slice/Ice/RouterF.ice
usr/share/ice/slice/Ice/ServantLocator.ice
usr/share/ice/slice/Ice/ServantLocatorF.ice
usr/share/ice/slice/Ice/SliceChecksumDict.ice
usr/share/ice/slice/Ice/ValueFactory.ice
usr/share/ice/slice/Ice/Version.ice
usr/share/ice/slice/IceBT/
usr/share/ice/slice/IceBT/ConnectionInfo.ice
usr/share/ice/slice/IceBT/EndpointInfo.ice
usr/share/ice/slice/IceBT/Types.ice
usr/share/ice/slice/IceBox/
usr/share/ice/slice/IceBox/IceBox.ice
usr/share/ice/slice/IceGrid/
usr/share/ice/slice/IceGrid/Admin.ice
usr/share/ice/slice/IceGrid/Descriptor.ice
usr/share/ice/slice/IceGrid/Exception.ice
usr/share/ice/slice/IceGrid/FileParser.ice
usr/share/ice/slice/IceGrid/PluginFacade.ice
usr/share/ice/slice/IceGrid/Registry.ice
usr/share/ice/slice/IceGrid/Session.ice
usr/share/ice/slice/IceGrid/UserAccountMapper.ice
usr/share/ice/slice/IceIAP/
usr/share/ice/slice/IceIAP/ConnectionInfo.ice
usr/share/ice/slice/IceIAP/EndpointInfo.ice
usr/share/ice/slice/IcePatch2/
usr/share/ice/slice/IcePatch2/FileInfo.ice
usr/share/ice/slice/IcePatch2/FileServer.ice
usr/share/ice/slice/IceSSL/
usr/share/ice/slice/IceSSL/ConnectionInfo.ice
usr/share/ice/slice/IceSSL/ConnectionInfoF.ice
usr/share/ice/slice/IceSSL/EndpointInfo.ice
usr/share/ice/slice/IceStorm/
usr/share/ice/slice/IceStorm/IceStorm.ice
usr/share/ice/slice/IceStorm/Metrics.ice
usr/share/ice/templates.xml
usr/share/licenses/
usr/share/licenses/zeroc-ice/
usr/share/licenses/zeroc-ice/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/glacier2router.1.gz
usr/share/man/man1/icebox.1.gz
usr/share/man/man1/iceboxadmin.1.gz
usr/share/man/man1/icebridge.1.gz
usr/share/man/man1/icegridadmin.1.gz
usr/share/man/man1/icegriddb.1.gz
usr/share/man/man1/icegridnode.1.gz
usr/share/man/man1/icegridregistry.1.gz
usr/share/man/man1/icepatch2calc.1.gz
usr/share/man/man1/icepatch2client.1.gz
usr/share/man/man1/icepatch2server.1.gz
usr/share/man/man1/icestormadmin.1.gz
usr/share/man/man1/icestormdb.1.gz
usr/share/man/man1/slice2cpp.1.gz
usr/share/man/man1/slice2cs.1.gz
usr/share/man/man1/slice2html.1.gz
usr/share/man/man1/slice2java.1.gz
usr/share/man/man1/slice2js.1.gz
usr/share/man/man1/slice2matlab.1.gz
usr/share/man/man1/slice2objc.1.gz
usr/share/man/man1/slice2php.1.gz
usr/share/man/man1/slice2rb.1.gz
usr/share/man/man1/slice2swift.1.gz
var/
var/lib/
var/lib/ice/
var/lib/ice/icegrid/
var/lib/ice/icegrid/node1/
var/lib/ice/icegrid/registry/
