# extracted/amavisd-new-2.12.2-2/files
%FILES%
etc/
etc/amavisd/
etc/amavisd/amavisd-custom.conf
etc/amavisd/amavisd.conf
etc/amavisd/amavisd.conf-default
etc/openldap/
etc/openldap/schema/
etc/openldap/schema/amavisd-new.ldif
etc/openldap/schema/amavisd-new.schema
usr/
usr/bin/
usr/bin/amavisd
usr/bin/amavisd-agent
usr/bin/amavisd-nanny
usr/bin/amavisd-release
usr/bin/amavisd-signer
usr/bin/amavisd-snmp-subagent
usr/bin/amavisd-snmp-subagent-zmq
usr/bin/p0f-analyzer.pl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/amavisd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/amavisd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/amavisd.conf
usr/share/
usr/share/doc/
usr/share/doc/amavis/
usr/share/doc/amavis/amavisd-new-courier.patch
usr/share/doc/amavis/amavisd-new-qmqpqq.patch
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/JpegTester.pm
var/
var/spool/
var/spool/amavis/
var/spool/amavis/db/
var/spool/amavis/tmp/
var/spool/amavis/var/
var/spool/amavis/virus/
