# extracted/prometheus-elasticsearch-exporter-1.7.0-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-elasticsearch-exporter
usr/
usr/bin/
usr/bin/prometheus-elasticsearch-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-elasticsearch-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-elasticsearch-exporter.conf
