# extracted/usbip-6.12-8/files
%FILES%
usr/
usr/bin/
usr/bin/usbip
usr/bin/usbipd
usr/include/
usr/include/usbip/
usr/include/usbip/list.h
usr/include/usbip/sysfs_utils.h
usr/include/usbip/usbip_common.h
usr/include/usbip/usbip_host_common.h
usr/include/usbip/usbip_host_driver.h
usr/include/usbip/vhci_driver.h
usr/lib/
usr/lib/libusbip.so
usr/lib/libusbip.so.0
usr/lib/libusbip.so.0.0.1
usr/lib/modules-load.d/
usr/lib/modules-load.d/usbip.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/usbipd.service
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/usbip.8.gz
usr/share/man/man8/usbipd.8.gz
