# extracted/pacman-offline-0.3.3-1/files
%FILES%
etc/
etc/pacman.d/
etc/pacman.d/offline.conf
usr/
usr/bin/
usr/bin/pacman-offline
usr/lib/
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/pacman-offline
usr/lib/systemd/system/
usr/lib/systemd/system/pacman-offline-prepare.service
usr/lib/systemd/system/pacman-offline-prepare.timer
usr/lib/systemd/system/pacman-offline-reboot.service
usr/lib/systemd/system/pacman-offline-reboot.timer
usr/lib/systemd/system/pacman-offline.service
usr/lib/systemd/system/system-update.target.wants/
usr/lib/systemd/system/system-update.target.wants/pacman-offline.service
usr/share/
usr/share/doc/
usr/share/doc/pacman-offline/
usr/share/doc/pacman-offline/README.html
usr/share/doc/pacman-offline/README.md
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/99-pacman-offline.hook
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/pacman-offline.rules
