# extracted/nextcloud-app-notify_push-0.7.0-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nextcloud-app-notify_push.service
usr/share/
usr/share/webapps/
usr/share/webapps/nextcloud/
usr/share/webapps/nextcloud/apps/
usr/share/webapps/nextcloud/apps/notify_push/
usr/share/webapps/nextcloud/apps/notify_push/appinfo/
usr/share/webapps/nextcloud/apps/notify_push/appinfo/info.xml
usr/share/webapps/nextcloud/apps/notify_push/appinfo/routes.php
usr/share/webapps/nextcloud/apps/notify_push/bin/
usr/share/webapps/nextcloud/apps/notify_push/bin/x86_64/
usr/share/webapps/nextcloud/apps/notify_push/bin/x86_64/notify_push
usr/share/webapps/nextcloud/apps/notify_push/lib/
usr/share/webapps/nextcloud/apps/notify_push/lib/AppInfo/
usr/share/webapps/nextcloud/apps/notify_push/lib/AppInfo/Application.php
usr/share/webapps/nextcloud/apps/notify_push/lib/BinaryFinder.php
usr/share/webapps/nextcloud/apps/notify_push/lib/CSPListener.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Capabilities.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/Log.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/Metrics.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/Reset.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/SelfTest.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Command/Setup.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Controller/
usr/share/webapps/nextcloud/apps/notify_push/lib/Controller/AuthController.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Controller/TestController.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Listener.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Migration/
usr/share/webapps/nextcloud/apps/notify_push/lib/Migration/Install.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Queue/
usr/share/webapps/nextcloud/apps/notify_push/lib/Queue/IQueue.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Queue/NullQueue.php
usr/share/webapps/nextcloud/apps/notify_push/lib/Queue/RedisQueue.php
usr/share/webapps/nextcloud/apps/notify_push/lib/SelfTest.php
usr/share/webapps/nextcloud/apps/notify_push/lib/SetupWizard.php
