# extracted/archivetools-4-1/files
%FILES%
etc/
etc/archive.conf
usr/
usr/bin/
usr/bin/archive
usr/bin/archive-cleaner
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/archive-hardlink.service
usr/lib/systemd/system/archive-hardlink.timer
usr/lib/systemd/system/archive.service
usr/lib/systemd/system/archive.timer
usr/lib/sysusers.d/
usr/lib/sysusers.d/archive.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/archive.conf
