# extracted/ollama-0.5.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/ollama
usr/lib/
usr/lib/ollama/
usr/lib/ollama/runners/
usr/lib/ollama/runners/cpu_avx/
usr/lib/ollama/runners/cpu_avx/ollama_llama_server
usr/lib/ollama/runners/cpu_avx2/
usr/lib/ollama/runners/cpu_avx2/ollama_llama_server
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ollama.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/ollama.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/ollama.conf
usr/share/
usr/share/licenses/
usr/share/licenses/ollama/
usr/share/licenses/ollama/LICENSE
usr/share/ollama
var/
var/lib/
var/lib/ollama/
