# extracted/phodav-3.0-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/spice-webdavd
usr/
usr/bin/
usr/bin/chezdav
usr/bin/spice-webdavd
usr/include/
usr/include/libphodav-3.0/
usr/include/libphodav-3.0/libphodav/
usr/include/libphodav-3.0/libphodav/phodav-server.h
usr/include/libphodav-3.0/libphodav/phodav-virtual-dir.h
usr/include/libphodav-3.0/libphodav/phodav.h
usr/lib/
usr/lib/libphodav-3.0.so
usr/lib/libphodav-3.0.so.0
usr/lib/libphodav-3.0.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libphodav-3.0.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/spice-webdavd.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/70-spice-webdavd.rules
usr/share/
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/phodav-3.0/
usr/share/gtk-doc/html/phodav-3.0/PhodavServer.html
usr/share/gtk-doc/html/phodav-3.0/PhodavVirtualDir.html
usr/share/gtk-doc/html/phodav-3.0/annotation-glossary.html
usr/share/gtk-doc/html/phodav-3.0/api-2.5.html
usr/share/gtk-doc/html/phodav-3.0/api-index-full.html
usr/share/gtk-doc/html/phodav-3.0/ch01.html
usr/share/gtk-doc/html/phodav-3.0/home.png
usr/share/gtk-doc/html/phodav-3.0/index.html
usr/share/gtk-doc/html/phodav-3.0/left-insensitive.png
usr/share/gtk-doc/html/phodav-3.0/left.png
usr/share/gtk-doc/html/phodav-3.0/object-tree.html
usr/share/gtk-doc/html/phodav-3.0/phodav-3.0.devhelp2
usr/share/gtk-doc/html/phodav-3.0/right-insensitive.png
usr/share/gtk-doc/html/phodav-3.0/right.png
usr/share/gtk-doc/html/phodav-3.0/style.css
usr/share/gtk-doc/html/phodav-3.0/up-insensitive.png
usr/share/gtk-doc/html/phodav-3.0/up.png
usr/share/locale/
usr/share/locale/bs/
usr/share/locale/bs/LC_MESSAGES/
usr/share/locale/bs/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/fur/
usr/share/locale/fur/LC_MESSAGES/
usr/share/locale/fur/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/sr@latin/
usr/share/locale/sr@latin/LC_MESSAGES/
usr/share/locale/sr@latin/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/phodav-3.0.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/phodav-3.0.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/chezdav.1.gz
