# extracted/prometheus-3.0.1-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus
etc/prometheus/
etc/prometheus/prometheus.yml
usr/
usr/bin/
usr/bin/prometheus
usr/bin/promtool
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus.conf
usr/share/
usr/share/doc/
usr/share/doc/prometheus/
usr/share/doc/prometheus/examples/
usr/share/doc/prometheus/examples/prometheus-agent.yml
usr/share/doc/prometheus/examples/prometheus-digitalocean.yml
usr/share/doc/prometheus/examples/prometheus-docker.yml
usr/share/doc/prometheus/examples/prometheus-dockerswarm.yml
usr/share/doc/prometheus/examples/prometheus-eureka.yml
usr/share/doc/prometheus/examples/prometheus-hetzner.yml
usr/share/doc/prometheus/examples/prometheus-ionos.yml
usr/share/doc/prometheus/examples/prometheus-kubernetes.yml
usr/share/doc/prometheus/examples/prometheus-linode.yml
usr/share/doc/prometheus/examples/prometheus-marathon.yml
usr/share/doc/prometheus/examples/prometheus-nomad.yml
usr/share/doc/prometheus/examples/prometheus-otlp.yml
usr/share/doc/prometheus/examples/prometheus-ovhcloud.yml
usr/share/doc/prometheus/examples/prometheus-puppetdb.yml
usr/share/doc/prometheus/examples/prometheus-scaleway.yml
usr/share/doc/prometheus/examples/prometheus-uyuni.yml
usr/share/doc/prometheus/examples/prometheus-vultr.yml
usr/share/doc/prometheus/examples/prometheus.yml
usr/share/prometheus/
usr/share/prometheus/web/
usr/share/prometheus/web/ui/
usr/share/prometheus/web/ui/static/
usr/share/prometheus/web/ui/static/mantine-ui/
usr/share/prometheus/web/ui/static/mantine-ui/assets/
usr/share/prometheus/web/ui/static/mantine-ui/assets/codicon-B_nZgZYP.ttf
usr/share/prometheus/web/ui/static/mantine-ui/assets/codicon-B_nZgZYP.ttf.gz
usr/share/prometheus/web/ui/static/mantine-ui/assets/index-1SrvQFrd.js
usr/share/prometheus/web/ui/static/mantine-ui/assets/index-1SrvQFrd.js.gz
usr/share/prometheus/web/ui/static/mantine-ui/assets/index-Bm6642QA.css
usr/share/prometheus/web/ui/static/mantine-ui/assets/index-Bm6642QA.css.gz
usr/share/prometheus/web/ui/static/mantine-ui/favicon.svg
usr/share/prometheus/web/ui/static/mantine-ui/favicon.svg.gz
usr/share/prometheus/web/ui/static/mantine-ui/index.html
usr/share/prometheus/web/ui/static/mantine-ui/index.html.gz
usr/share/prometheus/web/ui/static/react-app/
usr/share/prometheus/web/ui/static/react-app/asset-manifest.json
usr/share/prometheus/web/ui/static/react-app/asset-manifest.json.gz
usr/share/prometheus/web/ui/static/react-app/favicon.ico
usr/share/prometheus/web/ui/static/react-app/favicon.ico.gz
usr/share/prometheus/web/ui/static/react-app/index.html
usr/share/prometheus/web/ui/static/react-app/index.html.gz
usr/share/prometheus/web/ui/static/react-app/manifest.json
usr/share/prometheus/web/ui/static/react-app/manifest.json.gz
usr/share/prometheus/web/ui/static/react-app/static/
usr/share/prometheus/web/ui/static/react-app/static/css/
usr/share/prometheus/web/ui/static/react-app/static/css/main.a522f891.css
usr/share/prometheus/web/ui/static/react-app/static/css/main.a522f891.css.gz
usr/share/prometheus/web/ui/static/react-app/static/css/main.a522f891.css.map
usr/share/prometheus/web/ui/static/react-app/static/css/main.a522f891.css.map.gz
usr/share/prometheus/web/ui/static/react-app/static/js/
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js.LICENSE.txt
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js.LICENSE.txt.gz
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js.gz
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js.map
usr/share/prometheus/web/ui/static/react-app/static/js/main.46727cf5.js.map.gz
usr/share/prometheus/web/ui/static/react-app/static/media/
usr/share/prometheus/web/ui/static/react-app/static/media/codicon.b3726f0165bf67ac6849.ttf
usr/share/prometheus/web/ui/static/react-app/static/media/codicon.b3726f0165bf67ac6849.ttf.gz
usr/share/prometheus/web/ui/static/react-app/static/media/index.cd351d7c31d0d3fccf96.cjs
usr/share/prometheus/web/ui/static/react-app/static/media/index.cd351d7c31d0d3fccf96.cjs.gz
usr/share/prometheus/web/ui/static/react-app/static/media/prometheus_logo_grey.3cf697e5443028ca5e5255b93c7906c5.svg
usr/share/prometheus/web/ui/static/react-app/static/media/prometheus_logo_grey.3cf697e5443028ca5e5255b93c7906c5.svg.gz
var/
var/lib/
var/lib/prometheus/
