# extracted/openpgp-ca-restd-0.14.0-1/files
%FILES%
etc/
etc/openpgp-ca-restd/
etc/openpgp-ca-restd/env
usr/
usr/bin/
usr/bin/openpgp-ca-restd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openpgp-ca-restd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/openpgp-ca-restd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/openpgp-ca-restd.conf
