# extracted/gitea-1.22.6-1/files
%FILES%
etc/
etc/gitea/
etc/gitea/app.ini
usr/
usr/bin/
usr/bin/gitea
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gitea.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/gitea.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/gitea.conf
usr/share/
usr/share/licenses/
usr/share/licenses/gitea/
usr/share/licenses/gitea/LICENSE
