# extracted/jellyfin-server-10.10.3-1/files
%FILES%
etc/
etc/jellyfin/
etc/jellyfin/jellyfin.env
etc/jellyfin/logging.json
usr/
usr/bin/
usr/bin/jellyfin
usr/lib/
usr/lib/jellyfin/
usr/lib/jellyfin/ATL.dll
usr/lib/jellyfin/AsyncKeyedLock.dll
usr/lib/jellyfin/BDInfo.dll
usr/lib/jellyfin/BlurHashSharp.SkiaSharp.dll
usr/lib/jellyfin/BlurHashSharp.dll
usr/lib/jellyfin/CommandLine.dll
usr/lib/jellyfin/Diacritics.dll
usr/lib/jellyfin/DiscUtils.Core.dll
usr/lib/jellyfin/DiscUtils.Iso9660.dll
usr/lib/jellyfin/DiscUtils.Streams.dll
usr/lib/jellyfin/DiscUtils.Udf.dll
usr/lib/jellyfin/DotNet.Glob.dll
usr/lib/jellyfin/Emby.Naming.dll
usr/lib/jellyfin/Emby.Naming.xml
usr/lib/jellyfin/Emby.Photos.dll
usr/lib/jellyfin/Emby.Photos.xml
usr/lib/jellyfin/Emby.Server.Implementations.dll
usr/lib/jellyfin/Emby.Server.Implementations.xml
usr/lib/jellyfin/ExCSS.dll
usr/lib/jellyfin/HarfBuzzSharp.dll
usr/lib/jellyfin/ICU4N.Transliterator.dll
usr/lib/jellyfin/ICU4N.dll
usr/lib/jellyfin/J2N.dll
usr/lib/jellyfin/Jellyfin.Api.dll
usr/lib/jellyfin/Jellyfin.Api.xml
usr/lib/jellyfin/Jellyfin.Data.dll
usr/lib/jellyfin/Jellyfin.Data.xml
usr/lib/jellyfin/Jellyfin.Drawing.Skia.dll
usr/lib/jellyfin/Jellyfin.Drawing.Skia.xml
usr/lib/jellyfin/Jellyfin.Drawing.dll
usr/lib/jellyfin/Jellyfin.Drawing.xml
usr/lib/jellyfin/Jellyfin.Extensions.dll
usr/lib/jellyfin/Jellyfin.Extensions.xml
usr/lib/jellyfin/Jellyfin.LiveTv.dll
usr/lib/jellyfin/Jellyfin.LiveTv.xml
usr/lib/jellyfin/Jellyfin.MediaEncoding.Hls.dll
usr/lib/jellyfin/Jellyfin.MediaEncoding.Hls.xml
usr/lib/jellyfin/Jellyfin.MediaEncoding.Keyframes.dll
usr/lib/jellyfin/Jellyfin.MediaEncoding.Keyframes.xml
usr/lib/jellyfin/Jellyfin.Networking.dll
usr/lib/jellyfin/Jellyfin.Networking.xml
usr/lib/jellyfin/Jellyfin.Server.Implementations.dll
usr/lib/jellyfin/Jellyfin.Server.Implementations.xml
usr/lib/jellyfin/Jellyfin.XmlTv.dll
usr/lib/jellyfin/LrcParser.dll
usr/lib/jellyfin/MediaBrowser.Common.dll
usr/lib/jellyfin/MediaBrowser.Common.xml
usr/lib/jellyfin/MediaBrowser.Controller.dll
usr/lib/jellyfin/MediaBrowser.Controller.xml
usr/lib/jellyfin/MediaBrowser.LocalMetadata.dll
usr/lib/jellyfin/MediaBrowser.LocalMetadata.xml
usr/lib/jellyfin/MediaBrowser.MediaEncoding.dll
usr/lib/jellyfin/MediaBrowser.MediaEncoding.xml
usr/lib/jellyfin/MediaBrowser.Model.dll
usr/lib/jellyfin/MediaBrowser.Model.xml
usr/lib/jellyfin/MediaBrowser.Providers.dll
usr/lib/jellyfin/MediaBrowser.Providers.xml
usr/lib/jellyfin/MediaBrowser.XbmcMetadata.dll
usr/lib/jellyfin/MediaBrowser.XbmcMetadata.xml
usr/lib/jellyfin/MetaBrainz.Common.Json.dll
usr/lib/jellyfin/MetaBrainz.Common.dll
usr/lib/jellyfin/MetaBrainz.MusicBrainz.dll
usr/lib/jellyfin/Microsoft.AspNetCore.Authorization.dll
usr/lib/jellyfin/Microsoft.AspNetCore.Metadata.dll
usr/lib/jellyfin/Microsoft.Bcl.AsyncInterfaces.dll
usr/lib/jellyfin/Microsoft.Data.Sqlite.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Abstractions.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Relational.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.Sqlite.dll
usr/lib/jellyfin/Microsoft.EntityFrameworkCore.dll
usr/lib/jellyfin/Microsoft.Extensions.Caching.Memory.dll
usr/lib/jellyfin/Microsoft.Extensions.Configuration.Binder.dll
usr/lib/jellyfin/Microsoft.Extensions.Configuration.FileExtensions.dll
usr/lib/jellyfin/Microsoft.Extensions.Configuration.Json.dll
usr/lib/jellyfin/Microsoft.Extensions.DependencyInjection.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.DependencyInjection.dll
usr/lib/jellyfin/Microsoft.Extensions.DependencyModel.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.EntityFrameworkCore.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.HealthChecks.dll
usr/lib/jellyfin/Microsoft.Extensions.Diagnostics.dll
usr/lib/jellyfin/Microsoft.Extensions.Hosting.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.Http.dll
usr/lib/jellyfin/Microsoft.Extensions.Logging.Abstractions.dll
usr/lib/jellyfin/Microsoft.Extensions.Logging.dll
usr/lib/jellyfin/Microsoft.Extensions.Options.dll
usr/lib/jellyfin/Microsoft.OpenApi.dll
usr/lib/jellyfin/Microsoft.Win32.SystemEvents.dll
usr/lib/jellyfin/Mono.Nat.dll
usr/lib/jellyfin/NEbml.Core.dll
usr/lib/jellyfin/Newtonsoft.Json.dll
usr/lib/jellyfin/PlaylistsNET.dll
usr/lib/jellyfin/Prometheus.AspNetCore.dll
usr/lib/jellyfin/Prometheus.NetStandard.dll
usr/lib/jellyfin/Resources/
usr/lib/jellyfin/Resources/Configuration/
usr/lib/jellyfin/Resources/Configuration/logging.json
usr/lib/jellyfin/SQLitePCLRaw.batteries_v2.dll
usr/lib/jellyfin/SQLitePCLRaw.core.dll
usr/lib/jellyfin/SQLitePCLRaw.provider.e_sqlite3.dll
usr/lib/jellyfin/Serilog.AspNetCore.dll
usr/lib/jellyfin/Serilog.Enrichers.Thread.dll
usr/lib/jellyfin/Serilog.Extensions.Hosting.dll
usr/lib/jellyfin/Serilog.Extensions.Logging.dll
usr/lib/jellyfin/Serilog.Formatting.Compact.dll
usr/lib/jellyfin/Serilog.Settings.Configuration.dll
usr/lib/jellyfin/Serilog.Sinks.Async.dll
usr/lib/jellyfin/Serilog.Sinks.Console.dll
usr/lib/jellyfin/Serilog.Sinks.Debug.dll
usr/lib/jellyfin/Serilog.Sinks.File.dll
usr/lib/jellyfin/Serilog.Sinks.Graylog.Core.dll
usr/lib/jellyfin/Serilog.Sinks.Graylog.dll
usr/lib/jellyfin/Serilog.dll
usr/lib/jellyfin/ShimSkiaSharp.dll
usr/lib/jellyfin/SkiaSharp.HarfBuzz.dll
usr/lib/jellyfin/SkiaSharp.dll
usr/lib/jellyfin/Svg.Custom.dll
usr/lib/jellyfin/Svg.Model.dll
usr/lib/jellyfin/Svg.Skia.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.ReDoc.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.Swagger.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.SwaggerGen.dll
usr/lib/jellyfin/Swashbuckle.AspNetCore.SwaggerUI.dll
usr/lib/jellyfin/System.Drawing.Common.dll
usr/lib/jellyfin/System.Linq.Async.dll
usr/lib/jellyfin/TMDbLib.dll
usr/lib/jellyfin/TagLibSharp.dll
usr/lib/jellyfin/Ude.NetStandard.dll
usr/lib/jellyfin/UtfUnknown.dll
usr/lib/jellyfin/jellyfin
usr/lib/jellyfin/jellyfin.deps.json
usr/lib/jellyfin/jellyfin.dll
usr/lib/jellyfin/jellyfin.runtimeconfig.json
usr/lib/jellyfin/jellyfin.xml
usr/lib/jellyfin/libHarfBuzzSharp.so
usr/lib/jellyfin/libSkiaSharp.so
usr/lib/jellyfin/libe_sqlite3.so
usr/lib/jellyfin/libse.dll
usr/lib/jellyfin/prometheus-net.DotNetRuntime.dll
usr/lib/jellyfin/wwwroot/
usr/lib/jellyfin/wwwroot/api-docs/
usr/lib/jellyfin/wwwroot/api-docs/banner-dark.svg
usr/lib/jellyfin/wwwroot/api-docs/redoc/
usr/lib/jellyfin/wwwroot/api-docs/redoc/custom.css
usr/lib/jellyfin/wwwroot/api-docs/swagger/
usr/lib/jellyfin/wwwroot/api-docs/swagger/custom.css
usr/lib/jellyfin/zlib.net.dll
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/jellyfin.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/jellyfin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/jellyfin.conf
