# extracted/ferm-2.7-2/files
%FILES%
etc/
etc/ferm.conf
usr/
usr/bin/
usr/bin/ferm
usr/bin/import-ferm
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ferm.service
usr/share/
usr/share/doc/
usr/share/doc/ferm/
usr/share/doc/ferm/AUTHORS
usr/share/doc/ferm/COPYING
usr/share/doc/ferm/NEWS
usr/share/doc/ferm/README.rst
usr/share/doc/ferm/TODO
usr/share/doc/ferm/examples/
usr/share/doc/ferm/examples/antiddos.ferm
usr/share/doc/ferm/examples/arptables.ferm
usr/share/doc/ferm/examples/dmz_router.ferm
usr/share/doc/ferm/examples/dsl_router.ferm
usr/share/doc/ferm/examples/ebtables.ferm
usr/share/doc/ferm/examples/fileserver.ferm
usr/share/doc/ferm/examples/ipv6.ferm
usr/share/doc/ferm/examples/mailserver.ferm
usr/share/doc/ferm/examples/resolve.ferm
usr/share/doc/ferm/examples/webserver.ferm
usr/share/doc/ferm/examples/workstation.ferm
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ferm.1.gz
usr/share/man/man1/import-ferm.1.gz
