# extracted/rspamd-3.10.0-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/rspamd
etc/rspamd/
etc/rspamd/actions.conf
etc/rspamd/cgp.inc
etc/rspamd/common.conf
etc/rspamd/composites.conf
etc/rspamd/groups.conf
etc/rspamd/lang_detection.inc
etc/rspamd/logging.inc
etc/rspamd/maps.d/
etc/rspamd/maps.d/dmarc_whitelist.inc
etc/rspamd/maps.d/exe_clickbait.inc
etc/rspamd/maps.d/maillist.inc
etc/rspamd/maps.d/mid.inc
etc/rspamd/maps.d/mime_types.inc
etc/rspamd/maps.d/redirectors.inc
etc/rspamd/maps.d/spf_dkim_whitelist.inc
etc/rspamd/maps.d/surbl-whitelist.inc
etc/rspamd/metrics.conf
etc/rspamd/modules.conf
etc/rspamd/modules.d/
etc/rspamd/modules.d/antivirus.conf
etc/rspamd/modules.d/arc.conf
etc/rspamd/modules.d/asn.conf
etc/rspamd/modules.d/aws_s3.conf
etc/rspamd/modules.d/bayes_expiry.conf
etc/rspamd/modules.d/bimi.conf
etc/rspamd/modules.d/chartable.conf
etc/rspamd/modules.d/clickhouse.conf
etc/rspamd/modules.d/dcc.conf
etc/rspamd/modules.d/dkim.conf
etc/rspamd/modules.d/dkim_signing.conf
etc/rspamd/modules.d/dmarc.conf
etc/rspamd/modules.d/elastic.conf
etc/rspamd/modules.d/emails.conf
etc/rspamd/modules.d/external_relay.conf
etc/rspamd/modules.d/external_services.conf
etc/rspamd/modules.d/force_actions.conf
etc/rspamd/modules.d/forged_recipients.conf
etc/rspamd/modules.d/fuzzy_check.conf
etc/rspamd/modules.d/gpt.conf
etc/rspamd/modules.d/greylist.conf
etc/rspamd/modules.d/hfilter.conf
etc/rspamd/modules.d/history_redis.conf
etc/rspamd/modules.d/http_headers.conf
etc/rspamd/modules.d/known_senders.conf
etc/rspamd/modules.d/maillist.conf
etc/rspamd/modules.d/metadata_exporter.conf
etc/rspamd/modules.d/metric_exporter.conf
etc/rspamd/modules.d/mid.conf
etc/rspamd/modules.d/milter_headers.conf
etc/rspamd/modules.d/mime_types.conf
etc/rspamd/modules.d/multimap.conf
etc/rspamd/modules.d/mx_check.conf
etc/rspamd/modules.d/neural.conf
etc/rspamd/modules.d/once_received.conf
etc/rspamd/modules.d/p0f.conf
etc/rspamd/modules.d/phishing.conf
etc/rspamd/modules.d/ratelimit.conf
etc/rspamd/modules.d/rbl.conf
etc/rspamd/modules.d/redis.conf
etc/rspamd/modules.d/regexp.conf
etc/rspamd/modules.d/replies.conf
etc/rspamd/modules.d/reputation.conf
etc/rspamd/modules.d/rspamd_update.conf
etc/rspamd/modules.d/spamassassin.conf
etc/rspamd/modules.d/spamtrap.conf
etc/rspamd/modules.d/spf.conf
etc/rspamd/modules.d/surbl.conf
etc/rspamd/modules.d/trie.conf
etc/rspamd/modules.d/url_redirector.conf
etc/rspamd/modules.d/whitelist.conf
etc/rspamd/options.inc
etc/rspamd/rspamd.conf
etc/rspamd/scores.d/
etc/rspamd/scores.d/content_group.conf
etc/rspamd/scores.d/fuzzy_group.conf
etc/rspamd/scores.d/headers_group.conf
etc/rspamd/scores.d/hfilter_group.conf
etc/rspamd/scores.d/mime_types_group.conf
etc/rspamd/scores.d/mua_group.conf
etc/rspamd/scores.d/phishing_group.conf
etc/rspamd/scores.d/policies_group.conf
etc/rspamd/scores.d/rbl_group.conf
etc/rspamd/scores.d/statistics_group.conf
etc/rspamd/scores.d/subject_group.conf
etc/rspamd/scores.d/surbl_group.conf
etc/rspamd/scores.d/whitelist_group.conf
etc/rspamd/settings.conf
etc/rspamd/statistic.conf
etc/rspamd/worker-controller.inc
etc/rspamd/worker-fuzzy.inc
etc/rspamd/worker-normal.inc
etc/rspamd/worker-proxy.inc
usr/
usr/bin/
usr/bin/rspamadm
usr/bin/rspamadm-3.10.0
usr/bin/rspamc
usr/bin/rspamc-3.10.0
usr/bin/rspamd
usr/bin/rspamd-3.10.0
usr/bin/rspamd_stats
usr/lib/
usr/lib/rspamd/
usr/lib/rspamd/librspamd-actrie.so
usr/lib/rspamd/librspamd-ev.so
usr/lib/rspamd/librspamd-kann.so
usr/lib/rspamd/librspamd-replxx.so
usr/lib/rspamd/librspamd-server.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rspamd-dmarc-report.service
usr/lib/systemd/system/rspamd-dmarc-report.timer
usr/lib/systemd/system/rspamd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/rspamd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rspamd.conf
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/rspamd-hyperscan.hook
usr/share/licenses/
usr/share/licenses/rspamd/
usr/share/licenses/rspamd/LICENSE.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/rspamadm.1.gz
usr/share/man/man1/rspamc.1.gz
usr/share/man/man8/
usr/share/man/man8/rspamd.8.gz
usr/share/rspamd/
usr/share/rspamd/effective_tld_names.dat
usr/share/rspamd/elastic/
usr/share/rspamd/elastic/kibana.json
usr/share/rspamd/elastic/rspamd_template.json
usr/share/rspamd/languages/
usr/share/rspamd/languages/af.json
usr/share/rspamd/languages/an.json
usr/share/rspamd/languages/ar.json
usr/share/rspamd/languages/bg.json
usr/share/rspamd/languages/bn.json
usr/share/rspamd/languages/br.json
usr/share/rspamd/languages/cs.json
usr/share/rspamd/languages/cy.json
usr/share/rspamd/languages/da.json
usr/share/rspamd/languages/de.json
usr/share/rspamd/languages/en.json
usr/share/rspamd/languages/es.json
usr/share/rspamd/languages/et.json
usr/share/rspamd/languages/eu.json
usr/share/rspamd/languages/fa.json
usr/share/rspamd/languages/fi.json
usr/share/rspamd/languages/fr.json
usr/share/rspamd/languages/ga.json
usr/share/rspamd/languages/hi.json
usr/share/rspamd/languages/hr.json
usr/share/rspamd/languages/hu.json
usr/share/rspamd/languages/id.json
usr/share/rspamd/languages/is.json
usr/share/rspamd/languages/it.json
usr/share/rspamd/languages/lt.json
usr/share/rspamd/languages/lv.json
usr/share/rspamd/languages/mr.json
usr/share/rspamd/languages/ms.json
usr/share/rspamd/languages/ne.json
usr/share/rspamd/languages/nl.json
usr/share/rspamd/languages/no.json
usr/share/rspamd/languages/pa.json
usr/share/rspamd/languages/pl.json
usr/share/rspamd/languages/pt.json
usr/share/rspamd/languages/ro.json
usr/share/rspamd/languages/ru.json
usr/share/rspamd/languages/sl.json
usr/share/rspamd/languages/so.json
usr/share/rspamd/languages/sq.json
usr/share/rspamd/languages/sr.json
usr/share/rspamd/languages/stop_words
usr/share/rspamd/languages/sv.json
usr/share/rspamd/languages/sw.json
usr/share/rspamd/languages/tr.json
usr/share/rspamd/languages/uk.json
usr/share/rspamd/languages/ur.json
usr/share/rspamd/languages/vi.json
usr/share/rspamd/lualib/
usr/share/rspamd/lualib/ansicolors.lua
usr/share/rspamd/lualib/argparse.lua
usr/share/rspamd/lualib/fun.lua
usr/share/rspamd/lualib/global_functions.lua
usr/share/rspamd/lualib/lpegre.lua
usr/share/rspamd/lualib/lua_auth_results.lua
usr/share/rspamd/lualib/lua_aws.lua
usr/share/rspamd/lualib/lua_bayes_learn.lua
usr/share/rspamd/lualib/lua_bayes_redis.lua
usr/share/rspamd/lualib/lua_cfg_transform.lua
usr/share/rspamd/lualib/lua_cfg_utils.lua
usr/share/rspamd/lualib/lua_clickhouse.lua
usr/share/rspamd/lualib/lua_content/
usr/share/rspamd/lualib/lua_content/ical.lua
usr/share/rspamd/lualib/lua_content/init.lua
usr/share/rspamd/lualib/lua_content/pdf.lua
usr/share/rspamd/lualib/lua_content/vcard.lua
usr/share/rspamd/lualib/lua_dkim_tools.lua
usr/share/rspamd/lualib/lua_ffi/
usr/share/rspamd/lualib/lua_ffi/common.lua
usr/share/rspamd/lualib/lua_ffi/dkim.lua
usr/share/rspamd/lualib/lua_ffi/init.lua
usr/share/rspamd/lualib/lua_ffi/linalg.lua
usr/share/rspamd/lualib/lua_ffi/spf.lua
usr/share/rspamd/lualib/lua_fuzzy.lua
usr/share/rspamd/lualib/lua_lexer.lua
usr/share/rspamd/lualib/lua_magic/
usr/share/rspamd/lualib/lua_magic/heuristics.lua
usr/share/rspamd/lualib/lua_magic/init.lua
usr/share/rspamd/lualib/lua_magic/patterns.lua
usr/share/rspamd/lualib/lua_magic/types.lua
usr/share/rspamd/lualib/lua_maps.lua
usr/share/rspamd/lualib/lua_maps_expressions.lua
usr/share/rspamd/lualib/lua_meta.lua
usr/share/rspamd/lualib/lua_mime.lua
usr/share/rspamd/lualib/lua_mime_types.lua
usr/share/rspamd/lualib/lua_redis.lua
usr/share/rspamd/lualib/lua_scanners/
usr/share/rspamd/lualib/lua_scanners/avast.lua
usr/share/rspamd/lualib/lua_scanners/clamav.lua
usr/share/rspamd/lualib/lua_scanners/cloudmark.lua
usr/share/rspamd/lualib/lua_scanners/common.lua
usr/share/rspamd/lualib/lua_scanners/dcc.lua
usr/share/rspamd/lualib/lua_scanners/fprot.lua
usr/share/rspamd/lualib/lua_scanners/icap.lua
usr/share/rspamd/lualib/lua_scanners/init.lua
usr/share/rspamd/lualib/lua_scanners/kaspersky_av.lua
usr/share/rspamd/lualib/lua_scanners/kaspersky_se.lua
usr/share/rspamd/lualib/lua_scanners/oletools.lua
usr/share/rspamd/lualib/lua_scanners/p0f.lua
usr/share/rspamd/lualib/lua_scanners/pyzor.lua
usr/share/rspamd/lualib/lua_scanners/razor.lua
usr/share/rspamd/lualib/lua_scanners/savapi.lua
usr/share/rspamd/lualib/lua_scanners/sophos.lua
usr/share/rspamd/lualib/lua_scanners/spamassassin.lua
usr/share/rspamd/lualib/lua_scanners/vadesecure.lua
usr/share/rspamd/lualib/lua_scanners/virustotal.lua
usr/share/rspamd/lualib/lua_selectors/
usr/share/rspamd/lualib/lua_selectors/common.lua
usr/share/rspamd/lualib/lua_selectors/extractors.lua
usr/share/rspamd/lualib/lua_selectors/init.lua
usr/share/rspamd/lualib/lua_selectors/maps.lua
usr/share/rspamd/lualib/lua_selectors/transforms.lua
usr/share/rspamd/lualib/lua_settings.lua
usr/share/rspamd/lualib/lua_smtp.lua
usr/share/rspamd/lualib/lua_stat.lua
usr/share/rspamd/lualib/lua_tcp_sync.lua
usr/share/rspamd/lualib/lua_urls_compose.lua
usr/share/rspamd/lualib/lua_util.lua
usr/share/rspamd/lualib/lua_verdict.lua
usr/share/rspamd/lualib/lupa.lua
usr/share/rspamd/lualib/plugins/
usr/share/rspamd/lualib/plugins/dmarc.lua
usr/share/rspamd/lualib/plugins/neural.lua
usr/share/rspamd/lualib/plugins/ratelimit.lua
usr/share/rspamd/lualib/plugins/rbl.lua
usr/share/rspamd/lualib/plugins_stats.lua
usr/share/rspamd/lualib/redis_scripts/
usr/share/rspamd/lualib/redis_scripts/bayes_cache_check.lua
usr/share/rspamd/lualib/redis_scripts/bayes_cache_learn.lua
usr/share/rspamd/lualib/redis_scripts/bayes_classify.lua
usr/share/rspamd/lualib/redis_scripts/bayes_learn.lua
usr/share/rspamd/lualib/redis_scripts/bayes_stat.lua
usr/share/rspamd/lualib/redis_scripts/neural_maybe_invalidate.lua
usr/share/rspamd/lualib/redis_scripts/neural_maybe_lock.lua
usr/share/rspamd/lualib/redis_scripts/neural_save_unlock.lua
usr/share/rspamd/lualib/redis_scripts/neural_train_size.lua
usr/share/rspamd/lualib/redis_scripts/ratelimit_check.lua
usr/share/rspamd/lualib/redis_scripts/ratelimit_cleanup_pending.lua
usr/share/rspamd/lualib/redis_scripts/ratelimit_update.lua
usr/share/rspamd/lualib/rspamadm/
usr/share/rspamd/lualib/rspamadm/classifier_test.lua
usr/share/rspamd/lualib/rspamadm/clickhouse.lua
usr/share/rspamd/lualib/rspamadm/configgraph.lua
usr/share/rspamd/lualib/rspamadm/confighelp.lua
usr/share/rspamd/lualib/rspamadm/configwizard.lua
usr/share/rspamd/lualib/rspamadm/cookie.lua
usr/share/rspamd/lualib/rspamadm/corpus_test.lua
usr/share/rspamd/lualib/rspamadm/dkim_keygen.lua
usr/share/rspamd/lualib/rspamadm/dmarc_report.lua
usr/share/rspamd/lualib/rspamadm/dns_tool.lua
usr/share/rspamd/lualib/rspamadm/fuzzy_convert.lua
usr/share/rspamd/lualib/rspamadm/fuzzy_ping.lua
usr/share/rspamd/lualib/rspamadm/fuzzy_stat.lua
usr/share/rspamd/lualib/rspamadm/grep.lua
usr/share/rspamd/lualib/rspamadm/keypair.lua
usr/share/rspamd/lualib/rspamadm/mime.lua
usr/share/rspamd/lualib/rspamadm/neural_test.lua
usr/share/rspamd/lualib/rspamadm/publicsuffix.lua
usr/share/rspamd/lualib/rspamadm/secretbox.lua
usr/share/rspamd/lualib/rspamadm/stat_convert.lua
usr/share/rspamd/lualib/rspamadm/statistics_dump.lua
usr/share/rspamd/lualib/rspamadm/template.lua
usr/share/rspamd/lualib/rspamadm/vault.lua
usr/share/rspamd/lualib/tableshape.lua
usr/share/rspamd/plugins/
usr/share/rspamd/plugins/antivirus.lua
usr/share/rspamd/plugins/arc.lua
usr/share/rspamd/plugins/asn.lua
usr/share/rspamd/plugins/aws_s3.lua
usr/share/rspamd/plugins/bayes_expiry.lua
usr/share/rspamd/plugins/bimi.lua
usr/share/rspamd/plugins/clickhouse.lua
usr/share/rspamd/plugins/clustering.lua
usr/share/rspamd/plugins/dcc.lua
usr/share/rspamd/plugins/dkim_signing.lua
usr/share/rspamd/plugins/dmarc.lua
usr/share/rspamd/plugins/dynamic_conf.lua
usr/share/rspamd/plugins/elastic.lua
usr/share/rspamd/plugins/emails.lua
usr/share/rspamd/plugins/external_relay.lua
usr/share/rspamd/plugins/external_services.lua
usr/share/rspamd/plugins/force_actions.lua
usr/share/rspamd/plugins/forged_recipients.lua
usr/share/rspamd/plugins/fuzzy_collect.lua
usr/share/rspamd/plugins/gpt.lua
usr/share/rspamd/plugins/greylist.lua
usr/share/rspamd/plugins/hfilter.lua
usr/share/rspamd/plugins/history_redis.lua
usr/share/rspamd/plugins/http_headers.lua
usr/share/rspamd/plugins/ip_score.lua
usr/share/rspamd/plugins/known_senders.lua
usr/share/rspamd/plugins/maillist.lua
usr/share/rspamd/plugins/maps_stats.lua
usr/share/rspamd/plugins/metadata_exporter.lua
usr/share/rspamd/plugins/metric_exporter.lua
usr/share/rspamd/plugins/mid.lua
usr/share/rspamd/plugins/milter_headers.lua
usr/share/rspamd/plugins/mime_types.lua
usr/share/rspamd/plugins/multimap.lua
usr/share/rspamd/plugins/mx_check.lua
usr/share/rspamd/plugins/neural.lua
usr/share/rspamd/plugins/once_received.lua
usr/share/rspamd/plugins/p0f.lua
usr/share/rspamd/plugins/phishing.lua
usr/share/rspamd/plugins/ratelimit.lua
usr/share/rspamd/plugins/rbl.lua
usr/share/rspamd/plugins/replies.lua
usr/share/rspamd/plugins/reputation.lua
usr/share/rspamd/plugins/rspamd_update.lua
usr/share/rspamd/plugins/settings.lua
usr/share/rspamd/plugins/spamassassin.lua
usr/share/rspamd/plugins/spamtrap.lua
usr/share/rspamd/plugins/spf.lua
usr/share/rspamd/plugins/trie.lua
usr/share/rspamd/plugins/url_redirector.lua
usr/share/rspamd/plugins/whitelist.lua
usr/share/rspamd/rules/
usr/share/rspamd/rules/archives.lua
usr/share/rspamd/rules/bitcoin.lua
usr/share/rspamd/rules/bounce.lua
usr/share/rspamd/rules/content.lua
usr/share/rspamd/rules/controller/
usr/share/rspamd/rules/controller/fuzzy.lua
usr/share/rspamd/rules/controller/init.lua
usr/share/rspamd/rules/controller/maps.lua
usr/share/rspamd/rules/controller/neural.lua
usr/share/rspamd/rules/controller/selectors.lua
usr/share/rspamd/rules/forwarding.lua
usr/share/rspamd/rules/headers_checks.lua
usr/share/rspamd/rules/html.lua
usr/share/rspamd/rules/mid.lua
usr/share/rspamd/rules/misc.lua
usr/share/rspamd/rules/parts.lua
usr/share/rspamd/rules/regexp/
usr/share/rspamd/rules/regexp/compromised_hosts.lua
usr/share/rspamd/rules/regexp/headers.lua
usr/share/rspamd/rules/regexp/misc.lua
usr/share/rspamd/rules/regexp/upstream_spam_filters.lua
usr/share/rspamd/rules/rspamd.lua
usr/share/rspamd/rules/subject_checks.lua
usr/share/rspamd/www/
usr/share/rspamd/www/README.md
usr/share/rspamd/www/apple-touch-icon.png
usr/share/rspamd/www/browserconfig.xml
usr/share/rspamd/www/css/
usr/share/rspamd/www/css/FooTable.Glyphicons.css
usr/share/rspamd/www/css/bootstrap.min.css
usr/share/rspamd/www/css/codejar-linenumbers.css
usr/share/rspamd/www/css/d3evolution.css
usr/share/rspamd/www/css/d3pie.css
usr/share/rspamd/www/css/font-glyphicons.css
usr/share/rspamd/www/css/footable.standalone.min.css
usr/share/rspamd/www/css/nprogress.css
usr/share/rspamd/www/css/prism.css
usr/share/rspamd/www/css/rspamd.css
usr/share/rspamd/www/css/svg-with-js.min.css
usr/share/rspamd/www/favicon-16x16.png
usr/share/rspamd/www/favicon-32x32.png
usr/share/rspamd/www/favicon.ico
usr/share/rspamd/www/fonts/
usr/share/rspamd/www/fonts/glyphicons-halflings-regular.ttf
usr/share/rspamd/www/fonts/glyphicons-halflings-regular.woff
usr/share/rspamd/www/fonts/glyphicons-halflings-regular.woff2
usr/share/rspamd/www/img/
usr/share/rspamd/www/img/asc.png
usr/share/rspamd/www/img/desc.png
usr/share/rspamd/www/img/drop-area.svg
usr/share/rspamd/www/img/rspamd_logo_navbar.png
usr/share/rspamd/www/index.html
usr/share/rspamd/www/js/
usr/share/rspamd/www/js/app/
usr/share/rspamd/www/js/app/common.js
usr/share/rspamd/www/js/app/config.js
usr/share/rspamd/www/js/app/graph.js
usr/share/rspamd/www/js/app/history.js
usr/share/rspamd/www/js/app/libft.js
usr/share/rspamd/www/js/app/rspamd.js
usr/share/rspamd/www/js/app/selectors.js
usr/share/rspamd/www/js/app/stats.js
usr/share/rspamd/www/js/app/symbols.js
usr/share/rspamd/www/js/app/upload.js
usr/share/rspamd/www/js/lib/
usr/share/rspamd/www/js/lib/bootstrap.bundle.min.js
usr/share/rspamd/www/js/lib/codejar-linenumbers.min.js
usr/share/rspamd/www/js/lib/codejar.min.js
usr/share/rspamd/www/js/lib/d3.min.js
usr/share/rspamd/www/js/lib/d3evolution.min.js
usr/share/rspamd/www/js/lib/d3pie.min.js
usr/share/rspamd/www/js/lib/fontawesome.min.js
usr/share/rspamd/www/js/lib/footable.min.js
usr/share/rspamd/www/js/lib/jquery-3.7.1.min.js
usr/share/rspamd/www/js/lib/jquery.stickytabs.min.js
usr/share/rspamd/www/js/lib/nprogress.min.js
usr/share/rspamd/www/js/lib/prism.js
usr/share/rspamd/www/js/lib/require.min.js
usr/share/rspamd/www/js/lib/solid.min.js
usr/share/rspamd/www/js/lib/visibility.min.js
usr/share/rspamd/www/js/main.js
usr/share/rspamd/www/mstile-150x150.png
usr/share/rspamd/www/safari-pinned-tab.svg
