# extracted/prometheus-varnish-exporter-1.6.1-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-varnish-exporter
usr/
usr/bin/
usr/bin/prometheus-varnish-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-varnish-exporter.service
usr/share/
usr/share/doc/
usr/share/doc/prometheus-varnish-exporter/
usr/share/doc/prometheus-varnish-exporter/README.md
usr/share/licenses/
usr/share/licenses/prometheus-varnish-exporter/
usr/share/licenses/prometheus-varnish-exporter/LICENSE
