# extracted/dovecot-fts-xapian-1.8-1/files
%FILES%
usr/
usr/lib/
usr/lib/dovecot/
usr/lib/dovecot/modules/
usr/lib/dovecot/modules/lib21_fts_xapian_plugin.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dovecot-fts-optimize.service
usr/lib/systemd/system/dovecot-fts-optimize.timer
usr/share/
usr/share/doc/
usr/share/doc/dovecot-fts-xapian/
usr/share/doc/dovecot-fts-xapian/README.md
