# extracted/postfix-3.9.1-1/files
%FILES%
etc/
etc/postfix/
etc/postfix/LICENSE
etc/postfix/TLS_LICENSE
etc/postfix/access
etc/postfix/aliases
etc/postfix/bounce.cf.default
etc/postfix/canonical
etc/postfix/dynamicmaps.cf
etc/postfix/dynamicmaps.cf.d/
etc/postfix/generic
etc/postfix/header_checks
etc/postfix/main.cf
etc/postfix/main.cf.default
etc/postfix/main.cf.proto
etc/postfix/master.cf
etc/postfix/master.cf.proto
etc/postfix/postfix-files
etc/postfix/postfix-files.d/
etc/postfix/relocated
etc/postfix/transport
etc/postfix/virtual
usr/
usr/bin/
usr/bin/mailq
usr/bin/newaliases
usr/bin/postalias
usr/bin/postcat
usr/bin/postconf
usr/bin/postdrop
usr/bin/postfix
usr/bin/postfix-collate.pl
usr/bin/postfix-tlstype.pl
usr/bin/postkick
usr/bin/postlock
usr/bin/postlog
usr/bin/postmap
usr/bin/postmulti
usr/bin/postqueue
usr/bin/postsuper
usr/bin/posttls-finger
usr/bin/qmqp-sink
usr/bin/qmqp-source
usr/bin/qshape
usr/bin/sendmail
usr/bin/smtp-sink
usr/bin/smtp-source
usr/lib/
usr/lib/postfix/
usr/lib/postfix/bin/
usr/lib/postfix/bin/anvil
usr/lib/postfix/bin/bounce
usr/lib/postfix/bin/cleanup
usr/lib/postfix/bin/discard
usr/lib/postfix/bin/dnsblog
usr/lib/postfix/bin/error
usr/lib/postfix/bin/flush
usr/lib/postfix/bin/lmtp
usr/lib/postfix/bin/local
usr/lib/postfix/bin/master
usr/lib/postfix/bin/nqmgr
usr/lib/postfix/bin/oqmgr
usr/lib/postfix/bin/pickup
usr/lib/postfix/bin/pipe
usr/lib/postfix/bin/post-install
usr/lib/postfix/bin/postfix-script
usr/lib/postfix/bin/postfix-tls-script
usr/lib/postfix/bin/postfix-wrapper
usr/lib/postfix/bin/postlogd
usr/lib/postfix/bin/postmulti-script
usr/lib/postfix/bin/postscreen
usr/lib/postfix/bin/proxymap
usr/lib/postfix/bin/qmgr
usr/lib/postfix/bin/qmqpd
usr/lib/postfix/bin/scache
usr/lib/postfix/bin/showq
usr/lib/postfix/bin/smtp
usr/lib/postfix/bin/smtpd
usr/lib/postfix/bin/spawn
usr/lib/postfix/bin/tlsmgr
usr/lib/postfix/bin/tlsproxy
usr/lib/postfix/bin/trivial-rewrite
usr/lib/postfix/bin/verify
usr/lib/postfix/bin/virtual
usr/lib/postfix/libpostfix-dns.so
usr/lib/postfix/libpostfix-global.so
usr/lib/postfix/libpostfix-master.so
usr/lib/postfix/libpostfix-tls.so
usr/lib/postfix/libpostfix-util.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/postfix.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/postfix.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/postfix.conf
usr/share/
usr/share/doc/
usr/share/doc/postfix/
usr/share/doc/postfix/AAAREADME
usr/share/doc/postfix/ADDRESS_CLASS_README
usr/share/doc/postfix/ADDRESS_REWRITING_README
usr/share/doc/postfix/ADDRESS_VERIFICATION_README
usr/share/doc/postfix/BACKSCATTER_README
usr/share/doc/postfix/BASIC_CONFIGURATION_README
usr/share/doc/postfix/BDAT_README
usr/share/doc/postfix/BUILTIN_FILTER_README
usr/share/doc/postfix/COMPATIBILITY_README
usr/share/doc/postfix/CONNECTION_CACHE_README
usr/share/doc/postfix/CONTENT_INSPECTION_README
usr/share/doc/postfix/DATABASE_README
usr/share/doc/postfix/DB_README
usr/share/doc/postfix/DEBUG_README
usr/share/doc/postfix/DEPRECATION_README
usr/share/doc/postfix/DSN_README
usr/share/doc/postfix/ETRN_README
usr/share/doc/postfix/FILTER_README
usr/share/doc/postfix/FORWARD_SECRECY_README
usr/share/doc/postfix/INSTALL
usr/share/doc/postfix/IPV6_README
usr/share/doc/postfix/LINUX_README
usr/share/doc/postfix/LOCAL_RECIPIENT_README
usr/share/doc/postfix/MAILDROP_README
usr/share/doc/postfix/MAILLOG_README
usr/share/doc/postfix/MEMCACHE_README
usr/share/doc/postfix/MILTER_README
usr/share/doc/postfix/MULTI_INSTANCE_README
usr/share/doc/postfix/NFS_README
usr/share/doc/postfix/OVERVIEW
usr/share/doc/postfix/PACKAGE_README
usr/share/doc/postfix/POSTSCREEN_3_5_README
usr/share/doc/postfix/POSTSCREEN_README
usr/share/doc/postfix/QSHAPE_README
usr/share/doc/postfix/RELEASE_NOTES
usr/share/doc/postfix/RESTRICTION_CLASS_README
usr/share/doc/postfix/SASL_README
usr/share/doc/postfix/SCHEDULER_README
usr/share/doc/postfix/SMTPD_ACCESS_README
usr/share/doc/postfix/SMTPD_POLICY_README
usr/share/doc/postfix/SMTPD_PROXY_README
usr/share/doc/postfix/SMTPUTF8_README
usr/share/doc/postfix/SOHO_README
usr/share/doc/postfix/STANDARD_CONFIGURATION_README
usr/share/doc/postfix/STRESS_README
usr/share/doc/postfix/TLS_LEGACY_README
usr/share/doc/postfix/TLS_README
usr/share/doc/postfix/TUNING_README
usr/share/doc/postfix/ULTRIX_README
usr/share/doc/postfix/UUCP_README
usr/share/doc/postfix/VERP_README
usr/share/doc/postfix/VIRTUAL_README
usr/share/doc/postfix/XCLIENT_README
usr/share/doc/postfix/XFORWARD_README
usr/share/doc/postfix/collate/
usr/share/doc/postfix/collate/README
usr/share/doc/postfix/collate/README.tlstype
usr/share/licenses/
usr/share/licenses/postfix/
usr/share/licenses/postfix/COPYRIGHT
usr/share/licenses/postfix/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mailq.1.gz
usr/share/man/man1/newaliases.1.gz
usr/share/man/man1/postalias.1.gz
usr/share/man/man1/postcat.1.gz
usr/share/man/man1/postconf.1.gz
usr/share/man/man1/postdrop.1.gz
usr/share/man/man1/postfix-tls.1.gz
usr/share/man/man1/postfix.1.gz
usr/share/man/man1/postkick.1.gz
usr/share/man/man1/postlock.1.gz
usr/share/man/man1/postlog.1.gz
usr/share/man/man1/postmap.1.gz
usr/share/man/man1/postmulti.1.gz
usr/share/man/man1/postqueue.1.gz
usr/share/man/man1/postsuper.1.gz
usr/share/man/man1/posttls-finger.1.gz
usr/share/man/man1/qmqp-sink.1.gz
usr/share/man/man1/qmqp-source.1.gz
usr/share/man/man1/qshape.1.gz
usr/share/man/man1/sendmail.1.gz
usr/share/man/man1/smtp-sink.1.gz
usr/share/man/man1/smtp-source.1.gz
usr/share/man/man5/
usr/share/man/man5/access.5.gz
usr/share/man/man5/aliases.5.gz
usr/share/man/man5/body_checks.5.gz
usr/share/man/man5/bounce.5.gz
usr/share/man/man5/canonical.5.gz
usr/share/man/man5/cidr_table.5.gz
usr/share/man/man5/generic.5.gz
usr/share/man/man5/header_checks.5.gz
usr/share/man/man5/master.5.gz
usr/share/man/man5/memcache_table.5.gz
usr/share/man/man5/nisplus_table.5.gz
usr/share/man/man5/postconf.5.gz
usr/share/man/man5/postfix-wrapper.5.gz
usr/share/man/man5/regexp_table.5.gz
usr/share/man/man5/relocated.5.gz
usr/share/man/man5/socketmap_table.5.gz
usr/share/man/man5/tcp_table.5.gz
usr/share/man/man5/transport.5.gz
usr/share/man/man5/virtual.5.gz
usr/share/man/man8/
usr/share/man/man8/anvil.8.gz
usr/share/man/man8/bounce.8.gz
usr/share/man/man8/cleanup.8.gz
usr/share/man/man8/defer.8.gz
usr/share/man/man8/discard.8.gz
usr/share/man/man8/dnsblog.8.gz
usr/share/man/man8/error.8.gz
usr/share/man/man8/flush.8.gz
usr/share/man/man8/lmtp.8.gz
usr/share/man/man8/local.8.gz
usr/share/man/man8/master.8.gz
usr/share/man/man8/oqmgr.8.gz
usr/share/man/man8/pickup.8.gz
usr/share/man/man8/pipe.8.gz
usr/share/man/man8/postlogd.8.gz
usr/share/man/man8/postscreen.8.gz
usr/share/man/man8/proxymap.8.gz
usr/share/man/man8/qmgr.8.gz
usr/share/man/man8/qmqpd.8.gz
usr/share/man/man8/scache.8.gz
usr/share/man/man8/showq.8.gz
usr/share/man/man8/smtp.8.gz
usr/share/man/man8/smtpd.8.gz
usr/share/man/man8/spawn.8.gz
usr/share/man/man8/tlsmgr.8.gz
usr/share/man/man8/tlsproxy.8.gz
usr/share/man/man8/trace.8.gz
usr/share/man/man8/trivial-rewrite.8.gz
usr/share/man/man8/verify.8.gz
usr/share/man/man8/virtual.8.gz
var/
var/lib/
var/lib/postfix/
var/spool/
var/spool/postfix/
var/spool/postfix/active/
var/spool/postfix/bounce/
var/spool/postfix/corrupt/
var/spool/postfix/defer/
var/spool/postfix/deferred/
var/spool/postfix/flush/
var/spool/postfix/hold/
var/spool/postfix/incoming/
var/spool/postfix/maildrop/
var/spool/postfix/pid/
var/spool/postfix/private/
var/spool/postfix/public/
var/spool/postfix/saved/
var/spool/postfix/trace/
