# extracted/exim-4.98-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/exim
etc/mail/
etc/mail/aliases
etc/mail/exim.conf
usr/
usr/bin/
usr/bin/exicyclog
usr/bin/exigrep
usr/bin/exim
usr/bin/exim_checkaccess
usr/bin/exim_dbmbuild
usr/bin/exim_dumpdb
usr/bin/exim_fixdb
usr/bin/exim_id_update
usr/bin/exim_lock
usr/bin/exim_msgdate
usr/bin/exim_tidydb
usr/bin/eximstats
usr/bin/exinext
usr/bin/exipick
usr/bin/exiqgrep
usr/bin/exiqsumm
usr/bin/exiwhat
usr/bin/mailq
usr/bin/newaliases
usr/bin/rmail
usr/bin/rsmtp
usr/bin/runq
usr/bin/sendmail
usr/lib/
usr/lib/sendmail
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/exim-submission.socket
usr/lib/systemd/system/exim-submission@.service
usr/lib/systemd/system/exim.service
usr/lib/systemd/system/exim.socket
usr/lib/systemd/system/exim@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/exim.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/exim.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/exim.8.gz
