# extracted/mympd-18.2.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/mympd
usr/bin/mympd-config
usr/bin/mympd-script
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mympd.service
usr/lib/systemd/user/
usr/lib/systemd/user/mympd.service
usr/share/
usr/share/doc/
usr/share/doc/mympd/
usr/share/doc/mympd/CHANGELOG.md
usr/share/doc/mympd/LICENSE.md
usr/share/doc/mympd/README.md
usr/share/doc/mympd/SECURITY.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mympd-config.1.gz
usr/share/man/man1/mympd-script.1.gz
usr/share/man/man1/mympd.1.gz
