# extracted/nginx-prometheus-exporter-0.11.0-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/nginx-prometheus-exporter
usr/
usr/bin/
usr/bin/nginx-prometheus-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nginx-prometheus-exporter.service
