# extracted/hound-0.7.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/hound
usr/bin/houndd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hound.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/hound.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/hound.conf
usr/share/
usr/share/doc/
usr/share/doc/hound/
usr/share/doc/hound/CONTRIBUTING.md
usr/share/doc/hound/README.md
usr/share/doc/hound/config-example.json
usr/share/doc/hound/config-options.md
usr/share/doc/hound/default-config.json
usr/share/licenses/
usr/share/licenses/hound/
usr/share/licenses/hound/LICENSE
