# extracted/stubby-0.4.3-2/files
%FILES%
etc/
etc/stubby/
etc/stubby/stubby.yml
usr/
usr/bin/
usr/bin/stubby
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/stubby.service
usr/share/
usr/share/doc/
usr/share/doc/stubby/
usr/share/doc/stubby/AUTHORS
usr/share/doc/stubby/COPYING
usr/share/doc/stubby/ChangeLog
usr/share/doc/stubby/NEWS
usr/share/doc/stubby/README.md
usr/share/licenses/
usr/share/licenses/stubby/
usr/share/licenses/stubby/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/stubby.1.gz
