# extracted/ulogd-2.0.8-2/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/ulogd
etc/ulogd.conf
usr/
usr/bin/
usr/bin/ulogd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ulogd.service
usr/lib/ulogd/
usr/lib/ulogd/ulogd_filter_HWHDR.so
usr/lib/ulogd/ulogd_filter_IFINDEX.so
usr/lib/ulogd/ulogd_filter_IP2BIN.so
usr/lib/ulogd/ulogd_filter_IP2HBIN.so
usr/lib/ulogd/ulogd_filter_IP2STR.so
usr/lib/ulogd/ulogd_filter_MARK.so
usr/lib/ulogd/ulogd_filter_PRINTFLOW.so
usr/lib/ulogd/ulogd_filter_PRINTPKT.so
usr/lib/ulogd/ulogd_filter_PWSNIFF.so
usr/lib/ulogd/ulogd_inpflow_NFACCT.so
usr/lib/ulogd/ulogd_inpflow_NFCT.so
usr/lib/ulogd/ulogd_inppkt_NFLOG.so
usr/lib/ulogd/ulogd_inppkt_ULOG.so
usr/lib/ulogd/ulogd_inppkt_UNIXSOCK.so
usr/lib/ulogd/ulogd_output_DBI.so
usr/lib/ulogd/ulogd_output_GPRINT.so
usr/lib/ulogd/ulogd_output_GRAPHITE.so
usr/lib/ulogd/ulogd_output_IPFIX.so
usr/lib/ulogd/ulogd_output_JSON.so
usr/lib/ulogd/ulogd_output_LOGEMU.so
usr/lib/ulogd/ulogd_output_MYSQL.so
usr/lib/ulogd/ulogd_output_NACCT.so
usr/lib/ulogd/ulogd_output_OPRINT.so
usr/lib/ulogd/ulogd_output_PCAP.so
usr/lib/ulogd/ulogd_output_PGSQL.so
usr/lib/ulogd/ulogd_output_SQLITE3.so
usr/lib/ulogd/ulogd_output_SYSLOG.so
usr/lib/ulogd/ulogd_output_XML.so
usr/lib/ulogd/ulogd_raw2packet_BASE.so
usr/share/
usr/share/doc/
usr/share/doc/ulogd/
usr/share/doc/ulogd/mysql-ulogd2-flat.sql
usr/share/doc/ulogd/mysql-ulogd2.sql
usr/share/doc/ulogd/pgsql-ulogd2-flat.sql
usr/share/doc/ulogd/pgsql-ulogd2.sql
usr/share/doc/ulogd/sqlite3.table
usr/share/doc/ulogd/ulogd.conf
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/ulogd.8.gz
