# extracted/dnssec-tools-2.2.3-13/files
%FILES%
etc/
etc/dnssec-tools/
etc/dnssec-tools/KEY-SAFE/
etc/dnssec-tools/KEY-SAFE/.pacmankeep
etc/dnssec-tools/dnssec-tools.conf
etc/dnssec-tools/dnsval.conf
etc/dnssec-tools/root.hints
etc/dnssec-tools/validator-testcases
usr/
usr/bin/
usr/bin/dt-danechk
usr/bin/dt-getaddr
usr/bin/dt-gethost
usr/bin/dt-getname
usr/bin/dt-getquery
usr/bin/dt-getrrset
usr/bin/dt-libval_check_conf
usr/bin/dt-validate
usr/bin/libval-config
usr/bin/vendor_perl/
usr/bin/vendor_perl/blinkenlights
usr/bin/vendor_perl/bubbles
usr/bin/vendor_perl/buildrealms
usr/bin/vendor_perl/check-zone-expiration
usr/bin/vendor_perl/cleanarch
usr/bin/vendor_perl/cleankrf
usr/bin/vendor_perl/convertar
usr/bin/vendor_perl/dnspktflow
usr/bin/vendor_perl/donuts
usr/bin/vendor_perl/donutsd
usr/bin/vendor_perl/drawvalmap
usr/bin/vendor_perl/dtck
usr/bin/vendor_perl/dtconf
usr/bin/vendor_perl/dtconfchk
usr/bin/vendor_perl/dtdefs
usr/bin/vendor_perl/dtinitconf
usr/bin/vendor_perl/dtrealms
usr/bin/vendor_perl/expchk
usr/bin/vendor_perl/fixkrf
usr/bin/vendor_perl/genkrf
usr/bin/vendor_perl/getdnskeys
usr/bin/vendor_perl/getds
usr/bin/vendor_perl/grandvizier
usr/bin/vendor_perl/keyarch
usr/bin/vendor_perl/keymod
usr/bin/vendor_perl/krfcheck
usr/bin/vendor_perl/lights
usr/bin/vendor_perl/lsdnssec
usr/bin/vendor_perl/lskrf
usr/bin/vendor_perl/lsrealm
usr/bin/vendor_perl/lsroll
usr/bin/vendor_perl/maketestzone
usr/bin/vendor_perl/mapper
usr/bin/vendor_perl/realmchk
usr/bin/vendor_perl/realmctl
usr/bin/vendor_perl/realminit
usr/bin/vendor_perl/realmset
usr/bin/vendor_perl/rollchk
usr/bin/vendor_perl/rollctl
usr/bin/vendor_perl/rollerd
usr/bin/vendor_perl/rollinit
usr/bin/vendor_perl/rolllog
usr/bin/vendor_perl/rollrec-editor
usr/bin/vendor_perl/rollset
usr/bin/vendor_perl/signset-editor
usr/bin/vendor_perl/tachk
usr/bin/vendor_perl/timetrans
usr/bin/vendor_perl/trustman
usr/bin/vendor_perl/zonesigner
usr/include/
usr/include/validator/
usr/include/validator/resolver.h
usr/include/validator/val_dane.h
usr/include/validator/val_errors.h
usr/include/validator/validator-compat.h
usr/include/validator/validator-config.h
usr/include/validator/validator.h
usr/lib/
usr/lib/libsres.so
usr/lib/libsres.so.15
usr/lib/libsres.so.15.0.0
usr/lib/libval-threads.so
usr/lib/libval-threads.so.15
usr/lib/libval-threads.so.15.0.0
usr/lib/libval_shim.so
usr/lib/libval_shim.so.15
usr/lib/libval_shim.so.15.0.0
usr/lib/perl5/
usr/lib/perl5/5.40/
usr/lib/perl5/5.40/core_perl/
usr/lib/perl5/5.40/vendor_perl/
usr/lib/perl5/5.40/vendor_perl/Net/
usr/lib/perl5/5.40/vendor_perl/Net/DNS/
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/BootStrap.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/QWPrimitives.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/conf.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/defaults.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/dnssectools.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/keyrec.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/realm.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/realmmgr.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/rolllog.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/rollmgr.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/rollrec.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/timetrans.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Tools/tooloptions.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/Validator.pm
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/defines.pl
usr/lib/perl5/5.40/vendor_perl/Net/DNS/SEC/examples.pl
usr/lib/perl5/5.40/vendor_perl/Net/DNS/ZoneFile/
usr/lib/perl5/5.40/vendor_perl/Net/DNS/ZoneFile/Fast.pm
usr/lib/perl5/5.40/vendor_perl/Net/addrinfo.pm
usr/lib/perl5/5.40/vendor_perl/auto/
usr/lib/perl5/5.40/vendor_perl/auto/Net/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Dnspktflow/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Dnspktflow/dnspktflow/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Donuts/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Donuts/donuts/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Drawvalmap/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Drawvalmap/drawvalmap/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Maketestzone/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Maketestzone/maketestzone/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Mapper/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/Mapper/mapper/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Tools/conf/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Validator/
usr/lib/perl5/5.40/vendor_perl/auto/Net/DNS/SEC/Validator/Validator.so
usr/lib/perl5/5.40/vendor_perl/auto/Net/addrinfo/
usr/lib/perl5/5.40/vendor_perl/auto/Net/addrinfo/addrinfo.so
usr/lib/perl5/5.40/vendor_perl/auto/convertar/
usr/lib/perl5/5.40/vendor_perl/auto/dnssec-tools/
usr/lib/perl5/5.40/vendor_perl/auto/etc/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/donutsd.service
usr/lib/systemd/system/rollerd.service
usr/share/
usr/share/dnssec-tools/
usr/share/dnssec-tools/donuts/
usr/share/dnssec-tools/donuts/rules/
usr/share/dnssec-tools/donuts/rules/check_nameservers.txt
usr/share/dnssec-tools/donuts/rules/dns.errors.txt
usr/share/dnssec-tools/donuts/rules/dnssec.rules.txt
usr/share/dnssec-tools/donuts/rules/nsec_check.rules.txt
usr/share/dnssec-tools/donuts/rules/parent_child.rules.txt
usr/share/dnssec-tools/donuts/rules/recommendations.rules.txt
usr/share/dnssec-tools/donuts/rules/status_reports.txt
usr/share/licenses/
usr/share/licenses/dnssec-tools/
usr/share/licenses/dnssec-tools/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/blinkenlights.1p.gz
usr/share/man/man1/bubbles.1p.gz
usr/share/man/man1/buildrealms.1p.gz
usr/share/man/man1/check-zone-expiration.1p.gz
usr/share/man/man1/cleanarch.1p.gz
usr/share/man/man1/cleankrf.1p.gz
usr/share/man/man1/convertar.1p.gz
usr/share/man/man1/dnspktflow.1p.gz
usr/share/man/man1/dnssec-tools.1.gz
usr/share/man/man1/donuts.1p.gz
usr/share/man/man1/donutsd.1p.gz
usr/share/man/man1/drawvalmap.1p.gz
usr/share/man/man1/dt-danechk.1.gz
usr/share/man/man1/dt-getaddr.1.gz
usr/share/man/man1/dt-gethost.1.gz
usr/share/man/man1/dt-getname.1.gz
usr/share/man/man1/dt-getquery.1.gz
usr/share/man/man1/dt-getrrset.1.gz
usr/share/man/man1/dt-libval_check_conf.1.gz
usr/share/man/man1/dt-validate.1.gz
usr/share/man/man1/dtck.1p.gz
usr/share/man/man1/dtconf.1p.gz
usr/share/man/man1/dtconfchk.1p.gz
usr/share/man/man1/dtdefs.1p.gz
usr/share/man/man1/dtinitconf.1p.gz
usr/share/man/man1/dtrealms.1p.gz
usr/share/man/man1/expchk.1p.gz
usr/share/man/man1/fixkrf.1p.gz
usr/share/man/man1/genkrf.1p.gz
usr/share/man/man1/getdnskeys.1p.gz
usr/share/man/man1/getds.1p.gz
usr/share/man/man1/grandvizier.1p.gz
usr/share/man/man1/keyarch.1p.gz
usr/share/man/man1/keymod.1p.gz
usr/share/man/man1/krfcheck.1p.gz
usr/share/man/man1/lights.1p.gz
usr/share/man/man1/lsdnssec.1p.gz
usr/share/man/man1/lskrf.1p.gz
usr/share/man/man1/lsrealm.1p.gz
usr/share/man/man1/lsroll.1p.gz
usr/share/man/man1/maketestzone.1p.gz
usr/share/man/man1/mapper.1p.gz
usr/share/man/man1/realmchk.1p.gz
usr/share/man/man1/realmctl.1p.gz
usr/share/man/man1/realminit.1p.gz
usr/share/man/man1/realmset.1p.gz
usr/share/man/man1/rollchk.1p.gz
usr/share/man/man1/rollctl.1p.gz
usr/share/man/man1/rollerd.1p.gz
usr/share/man/man1/rollinit.1p.gz
usr/share/man/man1/rolllog.1p.gz
usr/share/man/man1/rollrec-editor.1p.gz
usr/share/man/man1/rollset.1p.gz
usr/share/man/man1/signset-editor.1p.gz
usr/share/man/man1/tachk.1p.gz
usr/share/man/man1/timetrans.1p.gz
usr/share/man/man1/trustman.1p.gz
usr/share/man/man1/zonesigner.1p.gz
usr/share/man/man3/
usr/share/man/man3/Net::DNS::SEC::Tools::BootStrap.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::Donuts.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::Donuts::Rule.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::QWPrimitives.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Bind.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Csv.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Dns.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Dump.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Itar.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Libval.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Mf.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::TrustAnchor::Secspider.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::conf.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::defaults.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::dnssectools.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::keyrec.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::realm.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::realmmgr.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::rolllog.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::rollmgr.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::rollrec.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::timetrans.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Tools::tooloptions.3pm.gz
usr/share/man/man3/Net::DNS::SEC::Validator.3pm.gz
usr/share/man/man3/Net::DNS::SEC::examples.3pm.gz
usr/share/man/man3/Net::DNS::ZoneFile::Fast.3pm.gz
usr/share/man/man3/Net::addrinfo.3pm.gz
usr/share/man/man3/dnsval.conf.3.gz
usr/share/man/man3/dnsval_conf_get.3.gz
usr/share/man/man3/dnsval_conf_set.3.gz
usr/share/man/man3/libsres.3.gz
usr/share/man/man3/libval.3.gz
usr/share/man/man3/libval_async.3.gz
usr/share/man/man3/libval_shim.3.gz
usr/share/man/man3/p_ac_status.3.gz
usr/share/man/man3/p_dane_error.3.gz
usr/share/man/man3/p_val_status.3.gz
usr/share/man/man3/resolv_conf_get.3.gz
usr/share/man/man3/resolv_conf_set.3.gz
usr/share/man/man3/root_hints_get.3.gz
usr/share/man/man3/root_hints_set.3.gz
usr/share/man/man3/val_add_valpolicy.3.gz
usr/share/man/man3/val_async_cancel.3.gz
usr/share/man/man3/val_async_cancel_all.3.gz
usr/share/man/man3/val_async_check_wait.3.gz
usr/share/man/man3/val_async_select_info.3.gz
usr/share/man/man3/val_async_submit.3.gz
usr/share/man/man3/val_context_setqflags.3.gz
usr/share/man/man3/val_create_context.3.gz
usr/share/man/man3/val_dane_check.3.gz
usr/share/man/man3/val_dane_match.3.gz
usr/share/man/man3/val_dane_submit.3.gz
usr/share/man/man3/val_does_not_exist.3.gz
usr/share/man/man3/val_free_answer_chain.3.gz
usr/share/man/man3/val_free_context.3.gz
usr/share/man/man3/val_free_dane.3.gz
usr/share/man/man3/val_free_response.3.gz
usr/share/man/man3/val_free_result_chain.3.gz
usr/share/man/man3/val_freeaddrinfo.3.gz
usr/share/man/man3/val_get_rrset.3.gz
usr/share/man/man3/val_getaddrinfo.3.gz
usr/share/man/man3/val_getdaneinfo.3.gz
usr/share/man/man3/val_gethostbyaddr.3.gz
usr/share/man/man3/val_gethostbyaddr_r.3.gz
usr/share/man/man3/val_gethostbyname.3.gz
usr/share/man/man3/val_gethostbyname2.3.gz
usr/share/man/man3/val_gethostbyname2_r.3.gz
usr/share/man/man3/val_gethostbyname_r.3.gz
usr/share/man/man3/val_getnameinfo.3.gz
usr/share/man/man3/val_istrusted.3.gz
usr/share/man/man3/val_isvalidated.3.gz
usr/share/man/man3/val_res_query.3.gz
usr/share/man/man3/val_res_search.3.gz
usr/share/man/man3/val_resolve_and_check.3.gz
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/Net/
usr/share/perl5/vendor_perl/Net/DNS/
usr/share/perl5/vendor_perl/Net/DNS/SEC/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/HTML.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/JSON.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/Perl.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/Text.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/Text/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/Text/Wrapped.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Output/Format/XML.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/Donuts/Rule.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Bind.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Csv.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Dns.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Dump.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Itar.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Libval.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Makefile.PL
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Mf.pm
usr/share/perl5/vendor_perl/Net/DNS/SEC/Tools/TrustAnchor/Secspider.pm
