# extracted/stratisd-3.6.3-1/files
%FILES%
usr/
usr/bin/
usr/bin/stratis-dumpmetadata
usr/bin/stratis-min
usr/bin/stratis-predict-usage
usr/lib/
usr/lib/dracut/
usr/lib/dracut/modules.d/
usr/lib/dracut/modules.d/90stratis-clevis/
usr/lib/dracut/modules.d/90stratis-clevis/module-setup.sh
usr/lib/dracut/modules.d/90stratis-clevis/stratis-clevis-rootfs-setup
usr/lib/dracut/modules.d/90stratis/
usr/lib/dracut/modules.d/90stratis/61-stratisd.rules
usr/lib/dracut/modules.d/90stratis/module-setup.sh
usr/lib/dracut/modules.d/90stratis/stratis-rootfs-setup
usr/lib/dracut/modules.d/90stratis/stratisd-min.service
usr/lib/stratisd
usr/lib/stratisd-min
usr/lib/systemd/
usr/lib/systemd/stratis-fstab-setup
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/stratis-clevis-setup-generator
usr/lib/systemd/system-generators/stratis-setup-generator
usr/lib/systemd/system/
usr/lib/systemd/system/stratis-fstab-setup@.service
usr/lib/systemd/system/stratisd-min-postinitrd.service
usr/lib/systemd/system/stratisd.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/61-stratisd.rules
usr/lib/udev/stratis-base32-decode
usr/lib/udev/stratis-str-cmp
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/stratisd.conf
usr/share/licenses/
usr/share/licenses/stratisd/
usr/share/licenses/stratisd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/stratis-dumpmetadata.8.gz
usr/share/man/man8/stratisd.8.gz
