# extracted/boinc-nox-8.0.2-2/files
%FILES%
usr/
usr/bin/
usr/bin/boinc
usr/bin/boinc_client
usr/bin/boinccmd
usr/include/
usr/include/boinc/
usr/include/boinc/app_ipc.h
usr/include/boinc/average.h
usr/include/boinc/base64.h
usr/include/boinc/boinc_api.h
usr/include/boinc/boinc_opencl.h
usr/include/boinc/boinc_stdio.h
usr/include/boinc/cal_boinc.h
usr/include/boinc/cc_config.h
usr/include/boinc/cert_sig.h
usr/include/boinc/cl_boinc.h
usr/include/boinc/common_defs.h
usr/include/boinc/coproc.h
usr/include/boinc/crypt.h
usr/include/boinc/diagnostics.h
usr/include/boinc/error_numbers.h
usr/include/boinc/filesys.h
usr/include/boinc/graphics2.h
usr/include/boinc/gui_rpc_client.h
usr/include/boinc/gutil.h
usr/include/boinc/hostinfo.h
usr/include/boinc/md5.h
usr/include/boinc/md5_file.h
usr/include/boinc/mem_usage.h
usr/include/boinc/mfile.h
usr/include/boinc/miofile.h
usr/include/boinc/msg_log.h
usr/include/boinc/network.h
usr/include/boinc/notice.h
usr/include/boinc/opencl_boinc.h
usr/include/boinc/parse.h
usr/include/boinc/prefs.h
usr/include/boinc/proc_control.h
usr/include/boinc/procinfo.h
usr/include/boinc/project_specific_defines.h
usr/include/boinc/proxy_info.h
usr/include/boinc/sched_msgs.h
usr/include/boinc/stackwalker_imports.h
usr/include/boinc/str_replace.h
usr/include/boinc/str_util.h
usr/include/boinc/svn_version.h
usr/include/boinc/url.h
usr/include/boinc/util.h
usr/include/boinc/version.h
usr/lib/
usr/lib/libboinc.a
usr/lib/libboinc_api.so
usr/lib/libboinc_api.so.8
usr/lib/libboinc_api.so.8.0.2
usr/lib/libboinc_crypt.a
usr/lib/libboinc_graphics2.so
usr/lib/libboinc_graphics2.so.8
usr/lib/libboinc_graphics2.so.8.0.2
usr/lib/libboinc_opencl.so
usr/lib/libboinc_opencl.so.8
usr/lib/libboinc_opencl.so.8.0.2
usr/lib/pkgconfig/
usr/lib/pkgconfig/libboinc.pc
usr/lib/pkgconfig/libboinc_api.pc
usr/lib/pkgconfig/libboinc_crypt.pc
usr/lib/pkgconfig/libboinc_graphics2.pc
usr/lib/pkgconfig/libboinc_opencl.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/boinc-client.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/boinc.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/boinc
