# extracted/spamassassin-4.0.1-2/files
%FILES%
etc/
etc/mail/
etc/mail/spamassassin/
etc/mail/spamassassin/init.pre
etc/mail/spamassassin/local.cf
etc/mail/spamassassin/v310.pre
etc/mail/spamassassin/v312.pre
etc/mail/spamassassin/v320.pre
etc/mail/spamassassin/v330.pre
etc/mail/spamassassin/v340.pre
etc/mail/spamassassin/v341.pre
etc/mail/spamassassin/v342.pre
etc/mail/spamassassin/v343.pre
etc/mail/spamassassin/v400.pre
etc/mail/spamassassin/v401.pre
usr/
usr/bin/
usr/bin/vendor_perl/
usr/bin/vendor_perl/sa-awl
usr/bin/vendor_perl/sa-check_spamd
usr/bin/vendor_perl/sa-compile
usr/bin/vendor_perl/sa-learn
usr/bin/vendor_perl/sa-update
usr/bin/vendor_perl/spamassassin
usr/bin/vendor_perl/spamc
usr/bin/vendor_perl/spamd
usr/lib/
usr/lib/perl5/
usr/lib/perl5/5.40/
usr/lib/perl5/5.40/core_perl/
usr/lib/perl5/5.40/vendor_perl/
usr/lib/perl5/5.40/vendor_perl/auto/
usr/lib/perl5/5.40/vendor_perl/auto/Mail/
usr/lib/perl5/5.40/vendor_perl/auto/Mail/SpamAssassin/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/spamassassin.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/spamassassin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/spamassassin.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/sa-awl.1p.gz
usr/share/man/man1/sa-check_spamd.1p.gz
usr/share/man/man1/sa-compile.1p.gz
usr/share/man/man1/sa-learn.1p.gz
usr/share/man/man1/sa-update.1p.gz
usr/share/man/man1/spamassassin-run.1p.gz
usr/share/man/man1/spamassassin.1p.gz
usr/share/man/man1/spamc.1p.gz
usr/share/man/man1/spamd.1p.gz
usr/share/man/man3/
usr/share/man/man3/Mail::SpamAssassin.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::AICache.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::ArchiveIterator.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::AsyncLoop.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::AutoWelcomelist.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Bayes.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore::BDB.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore::MySQL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore::PgSQL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore::Redis.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::BayesStore::SQL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Client.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Conf.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Conf::LDAP.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Conf::Parser.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Conf::SQL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::DnsResolver.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::GeoDB.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Logger.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Logger::File.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Logger::Stderr.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Logger::Syslog.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Message.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Message::Metadata.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Message::Node.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::PerMsgLearner.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::PerMsgStatus.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::PersistentAddrList.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::ASN.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AWL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AccessDB.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AntiVirus.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AskDNS.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AuthRes.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::AutoLearnThreshold.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Bayes.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::BodyRuleBaseExtractor.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Check.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::DCC.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::DKIM.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::DMARC.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::DNSEval.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::DecodeShortURLs.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::ExtractText.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::FreeMail.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::FromNameSpoof.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::HashBL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::MIMEEval.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::MIMEHeader.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::OLEVBMacro.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::OneLineBodyRuleType.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::PDFInfo.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::PhishTag.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Phishing.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Pyzor.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Razor2.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::RelayCountry.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::ReplaceTags.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::ResourceLimits.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Reuse.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Rule2XSBody.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::SPF.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Shortcircuit.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::SpamCop.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::Test.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::TextCat.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::TxRep.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::URIDNSBL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::URIDetail.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::URILocalBL.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::VBounce.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Plugin::WelcomeListSubject.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::PluginHandler.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Pyzor.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Pyzor::Client.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Pyzor::Digest.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Pyzor::Digest::Pieces.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Pyzor::Digest::StripHtml.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::RegistryBoundaries.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::SQLBasedAddrList.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::SubProcBackChannel.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Timeout.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Util.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Util::DependencyInfo.3pm.gz
usr/share/man/man3/Mail::SpamAssassin::Util::Progress.3pm.gz
usr/share/man/man3/spamassassin-run.3pm.gz
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/Mail/
usr/share/perl5/vendor_perl/Mail/SpamAssassin.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/AICache.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/ArchiveIterator.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/AsyncLoop.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/AutoWelcomelist.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Bayes.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Bayes/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Bayes/CombineChi.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Bayes/CombineNaiveBayes.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/BDB.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/DBM.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/MySQL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/PgSQL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/Redis.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/SDBM.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/BayesStore/SQL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Client.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Conf.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Conf/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Conf/LDAP.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Conf/Parser.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Conf/SQL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Constants.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/DBBasedAddrList.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Dns.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/DnsResolver.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/GeoDB.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/HTML.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locales.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locker.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locker/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locker/Flock.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locker/UnixNFSSafe.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Locker/Win32.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Logger.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Logger/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Logger/File.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Logger/Stderr.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Logger/Syslog.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/MailingList.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message/Metadata.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message/Metadata/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message/Metadata/Received.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Message/Node.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/NetSet.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/PerMsgLearner.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/PerMsgStatus.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/PersistentAddrList.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/ASN.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AWL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AccessDB.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AntiVirus.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AskDNS.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AuthRes.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/AutoLearnThreshold.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Bayes.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/BodyEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/BodyRuleBaseExtractor.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Check.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/DCC.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/DKIM.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/DMARC.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/DNSEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/DecodeShortURLs.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/ExtractText.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/FreeMail.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/FromNameSpoof.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/HTMLEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/HTTPSMismatch.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/HashBL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/HeaderEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/ImageInfo.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/MIMEEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/MIMEHeader.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/OLEVBMacro.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/OneLineBodyRuleType.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/PDFInfo.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/PhishTag.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Phishing.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Pyzor.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Razor2.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/RelayCountry.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/RelayEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/ReplaceTags.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/ResourceLimits.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Reuse.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Rule2XSBody.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/SPF.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Shortcircuit.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/SpamCop.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/Test.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/TextCat.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/TxRep.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/URIDNSBL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/URIDetail.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/URIEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/URILocalBL.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/VBounce.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/WLBLEval.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Plugin/WelcomeListSubject.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/PluginHandler.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/Client.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/Digest.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/Digest/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/Digest/Pieces.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Pyzor/Digest/StripHtml.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/RegistryBoundaries.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Reporter.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/SQLBasedAddrList.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/SpamdForkScaling.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/SubProcBackChannel.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Timeout.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/DependencyInfo.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/Progress.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/ScopedTimer.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/TieOneStringHash.pm
usr/share/perl5/vendor_perl/Mail/SpamAssassin/Util/TinyRedis.pm
usr/share/spamassassin/
usr/share/spamassassin/languages
usr/share/spamassassin/sa-update-pubkey.txt
usr/share/spamassassin/user_prefs.template
var/
var/lib/
var/lib/spamassassin/
