# extracted/lib32-pcsclite-2.3.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/pcscd-32
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pcscd-32.service
usr/lib/systemd/system/pcscd-32.socket
usr/lib32/
usr/lib32/libpcsclite.so
usr/lib32/libpcsclite.so.1
usr/lib32/libpcsclite_real.so
usr/lib32/libpcsclite_real.so.1
usr/lib32/libpcscspy.so
usr/lib32/libpcscspy.so.0
usr/lib32/pcsc/
usr/lib32/pcsc/drivers/
usr/lib32/pkgconfig/
usr/lib32/pkgconfig/libpcsclite.pc
usr/share/
usr/share/licenses/
usr/share/licenses/lib32-pcsclite/
usr/share/licenses/lib32-pcsclite/LICENSE.0BSD
usr/share/licenses/lib32-pcsclite/LICENSE.BSD-2-Clause
usr/share/licenses/lib32-pcsclite/LICENSE.BSD-3-Clause
usr/share/licenses/lib32-pcsclite/LICENSE.GPL-3.0-or-later
