# extracted/glider-0.16.4-1/files
%FILES%
etc/
etc/glider/
etc/glider/glider.conf.example
usr/
usr/bin/
usr/bin/glider
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/glider@.service
