# extracted/soft-serve-0.8.1-1/files
%FILES%
etc/
etc/soft-serve.conf
usr/
usr/bin/
usr/bin/soft
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/soft-serve.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/soft-serve.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/soft-serve.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/soft
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/soft.fish
usr/share/licenses/
usr/share/licenses/soft-serve/
usr/share/licenses/soft-serve/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/soft.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_soft
