# extracted/privoxy-3.0.34-3/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/privoxy
etc/privoxy/
etc/privoxy/config
etc/privoxy/default.action
etc/privoxy/default.filter
etc/privoxy/match-all.action
etc/privoxy/regression-tests.action
etc/privoxy/templates/
etc/privoxy/templates/blocked
etc/privoxy/templates/cgi-error-404
etc/privoxy/templates/cgi-error-bad-param
etc/privoxy/templates/cgi-error-disabled
etc/privoxy/templates/cgi-error-file
etc/privoxy/templates/cgi-error-file-read-only
etc/privoxy/templates/cgi-error-modified
etc/privoxy/templates/cgi-error-parse
etc/privoxy/templates/cgi-style.css
etc/privoxy/templates/client-tags
etc/privoxy/templates/connect-failed
etc/privoxy/templates/connection-timeout
etc/privoxy/templates/default
etc/privoxy/templates/edit-actions-add-url-form
etc/privoxy/templates/edit-actions-for-url
etc/privoxy/templates/edit-actions-for-url-filter
etc/privoxy/templates/edit-actions-for-url-string-action
etc/privoxy/templates/edit-actions-list
etc/privoxy/templates/edit-actions-list-button
etc/privoxy/templates/edit-actions-list-section
etc/privoxy/templates/edit-actions-list-url
etc/privoxy/templates/edit-actions-remove-url-form
etc/privoxy/templates/edit-actions-url-form
etc/privoxy/templates/forwarding-failed
etc/privoxy/templates/mod-local-help
etc/privoxy/templates/mod-support-and-service
etc/privoxy/templates/mod-title
etc/privoxy/templates/mod-unstable-warning
etc/privoxy/templates/no-server-data
etc/privoxy/templates/no-such-domain
etc/privoxy/templates/show-request
etc/privoxy/templates/show-status
etc/privoxy/templates/show-status-file
etc/privoxy/templates/show-url-info
etc/privoxy/templates/toggle
etc/privoxy/templates/toggle-mini
etc/privoxy/templates/untrusted
etc/privoxy/templates/url-info-osd.xml
etc/privoxy/templates/wpad.dat
etc/privoxy/trust
etc/privoxy/user.action
etc/privoxy/user.filter
usr/
usr/bin/
usr/bin/privoxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/privoxy.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/privoxy.conf
usr/share/
usr/share/doc/
usr/share/doc/privoxy/
usr/share/doc/privoxy/AUTHORS
usr/share/doc/privoxy/ChangeLog
usr/share/doc/privoxy/README
usr/share/doc/privoxy/developer-manual/
usr/share/doc/privoxy/developer-manual/coding.html
usr/share/doc/privoxy/developer-manual/documentation.html
usr/share/doc/privoxy/developer-manual/git.html
usr/share/doc/privoxy/developer-manual/index.html
usr/share/doc/privoxy/developer-manual/introduction.html
usr/share/doc/privoxy/developer-manual/newrelease.html
usr/share/doc/privoxy/developer-manual/testing.html
usr/share/doc/privoxy/developer-manual/webserver-update.html
usr/share/doc/privoxy/faq/
usr/share/doc/privoxy/faq/configuration.html
usr/share/doc/privoxy/faq/contact.html
usr/share/doc/privoxy/faq/copyright.html
usr/share/doc/privoxy/faq/general.html
usr/share/doc/privoxy/faq/index.html
usr/share/doc/privoxy/faq/installation.html
usr/share/doc/privoxy/faq/misc.html
usr/share/doc/privoxy/faq/trouble.html
usr/share/doc/privoxy/index.html
usr/share/doc/privoxy/man-page/
usr/share/doc/privoxy/man-page/privoxy-man-page.html
usr/share/doc/privoxy/p_doc.css
usr/share/doc/privoxy/user-manual/
usr/share/doc/privoxy/user-manual/actions-file.html
usr/share/doc/privoxy/user-manual/appendix.html
usr/share/doc/privoxy/user-manual/config.html
usr/share/doc/privoxy/user-manual/configuration.html
usr/share/doc/privoxy/user-manual/contact.html
usr/share/doc/privoxy/user-manual/copyright.html
usr/share/doc/privoxy/user-manual/files-in-use.jpg
usr/share/doc/privoxy/user-manual/filter-file.html
usr/share/doc/privoxy/user-manual/index.html
usr/share/doc/privoxy/user-manual/installation.html
usr/share/doc/privoxy/user-manual/introduction.html
usr/share/doc/privoxy/user-manual/p_doc.css
usr/share/doc/privoxy/user-manual/proxy2.jpg
usr/share/doc/privoxy/user-manual/proxy_setup.jpg
usr/share/doc/privoxy/user-manual/quickstart.html
usr/share/doc/privoxy/user-manual/seealso.html
usr/share/doc/privoxy/user-manual/startup.html
usr/share/doc/privoxy/user-manual/templates.html
usr/share/doc/privoxy/user-manual/whatsnew.html
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/privoxy.8.gz
var/
var/log/
var/log/privoxy/
