# extracted/emptty-0.13.0-1/files
%FILES%
etc/
etc/emptty/
etc/emptty/conf
etc/pam.d/
etc/pam.d/emptty
usr/
usr/bin/
usr/bin/emptty
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/emptty.service
usr/share/
usr/share/doc/
usr/share/doc/emptty/
usr/share/doc/emptty/README.md
usr/share/doc/emptty/SAMPLES.md
usr/share/doc/emptty/screenshot.png
usr/share/licenses/
usr/share/licenses/emptty/
usr/share/licenses/emptty/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/emptty.1.gz
