# extracted/deepin-boot-maker-5.10.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/deepin-boot-maker
usr/lib/
usr/lib/deepin-daemon/
usr/lib/deepin-daemon/deepin-boot-maker-service
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/deepin-boot-maker.service
usr/share/
usr/share/applications/
usr/share/applications/deepin-boot-maker.desktop
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/com.deepin.bootmaker.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.deepin.bootmaker.conf
usr/share/deepin-boot-maker/
usr/share/deepin-boot-maker/translations/
usr/share/deepin-boot-maker/translations/deepin-boot-maker.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_af.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_am_ET.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ar.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ast.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_az.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_bg.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_bn.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_bo.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_br.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ca.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_cs.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_da.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_de.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_el.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_en_AU.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_eo.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_es.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_fa.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_fi.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_fr.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_gl_ES.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_he.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_hi_IN.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_hr.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_hu.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_hy.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_id.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_it.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ja.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_km_KH.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_kn_IN.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ko.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ku_IQ.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ky.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ky@Arab.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_lo.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_lt.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ml.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_mn.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ms.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_nb.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ne.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_nl.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_pam.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_pl.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_pt.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_pt_BR.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ro.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ru.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_sk.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_sl.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_sq.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_sr.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_sv.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ta.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_th.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_tr.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_ug.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_uk.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_vi.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_zh_CN.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_zh_HK.qm
usr/share/deepin-boot-maker/translations/deepin-boot-maker_zh_TW.qm
usr/share/deepin-manual/
usr/share/deepin-manual/manual-assets/
usr/share/deepin-manual/manual-assets/application/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/\345\220\257\345\212\250\347\233\230\345\210\266\344\275\234\345\267\245\345\205\267\345\220\204\347\211\210\346\234\254\345\267\256\345\274\202\345\217\212\345\217\230\346\233\264\350\256\260\345\275\225.txt
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/attention.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/close.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/deepin_boot_maker.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/deepin_boot_maker_96.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/deepin_launcher.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/icon_close.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/icon_menu.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/notes.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/common/tips.svg
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/aftersale_services.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/drive.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/failed.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/making.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/select_disk.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/selectfile.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/fig/success.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/en_US/p_boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/boot_making.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/failed.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/failed_01.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/select_disk.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/select_file.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/fig/success.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_CN/p_boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/boot_making.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/failed.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/failed_01.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/select_disk.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/select_file.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/fig/success.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_HK/p_boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/boot-maker.md
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/aftersale_services.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/fail.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/failure.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/making.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/select_disk.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/selectfile.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/fig/success.png
usr/share/deepin-manual/manual-assets/application/deepin-boot-maker/boot-maker/zh_TW/p_boot-maker.md
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/deepin-boot-maker.svg
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/com.deepin.bootmaker.policy
