# extracted/nullmailer-2.2-4/files
%FILES%
etc/
etc/nullmailer/
usr/
usr/bin/
usr/bin/mailq
usr/bin/nullmailer-dsn
usr/bin/nullmailer-inject
usr/bin/nullmailer-queue
usr/bin/nullmailer-send
usr/bin/nullmailer-smtpd
usr/bin/sendmail
usr/lib/
usr/lib/nullmailer/
usr/lib/nullmailer/nullmailer/
usr/lib/nullmailer/nullmailer/qmqp
usr/lib/nullmailer/nullmailer/smtp
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nullmailer.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nullmailer.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nullmailer.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nullmailer-dsn.1.gz
usr/share/man/man1/nullmailer-inject.1.gz
usr/share/man/man1/sendmail.1.gz
usr/share/man/man7/
usr/share/man/man7/nullmailer.7.gz
usr/share/man/man8/
usr/share/man/man8/nullmailer-queue.8.gz
usr/share/man/man8/nullmailer-send.8.gz
