# extracted/tempo-2.5.0-4/files
%FILES%
etc/
etc/tempo/
etc/tempo/config.yml
usr/
usr/bin/
usr/bin/tempo
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tempo.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/tempo.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/tempo.conf
