# extracted/dbmail-3.4.1-3/files
%FILES%
etc/
etc/dbmail.conf.sample
usr/
usr/bin/
usr/bin/dbmail-deliver
usr/bin/dbmail-export
usr/bin/dbmail-httpd
usr/bin/dbmail-imapd
usr/bin/dbmail-lmtpd
usr/bin/dbmail-pop3d
usr/bin/dbmail-sievecmd
usr/bin/dbmail-sieved
usr/bin/dbmail-timsieved
usr/bin/dbmail-users
usr/bin/dbmail-util
usr/lib/
usr/lib/dbmail/
usr/lib/dbmail/libauth_ldap.so
usr/lib/dbmail/libauth_ldap.so.0
usr/lib/dbmail/libauth_ldap.so.0.0.0
usr/lib/dbmail/libauth_sql.so
usr/lib/dbmail/libauth_sql.so.0
usr/lib/dbmail/libauth_sql.so.0.0.0
usr/lib/dbmail/libdbmail.so
usr/lib/dbmail/libdbmail.so.0
usr/lib/dbmail/libdbmail.so.0.0.0
usr/lib/dbmail/libsort_sieve.so
usr/lib/dbmail/libsort_sieve.so.0
usr/lib/dbmail/libsort_sieve.so.0.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dbmail-imapd.service
usr/lib/systemd/system/dbmail-lmtpd.service
usr/lib/systemd/system/dbmail-pop3d.service
usr/lib/systemd/system/dbmail-timsieved.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/dbmail.conf
usr/share/
usr/share/dbmail/
usr/share/dbmail/contrib/
usr/share/dbmail/contrib/autosubscribe/
usr/share/dbmail/contrib/autosubscribe/autosubscribe.sql
usr/share/dbmail/contrib/build-centos-rpm/
usr/share/dbmail/contrib/build-centos-rpm/dbmail.spec
usr/share/dbmail/contrib/build-centos-rpm/libsieve.spec
usr/share/dbmail/contrib/build-centos-rpm/libzdb.spec
usr/share/dbmail/contrib/build-deb-with-docker.sh
usr/share/dbmail/contrib/build-rpm/
usr/share/dbmail/contrib/build-rpm/RPM_README
usr/share/dbmail/contrib/build-rpm/dbmail.spec
usr/share/dbmail/contrib/build-rpm/libsieve.spec
usr/share/dbmail/contrib/dbmail.spec
usr/share/dbmail/contrib/dbmail.supp
usr/share/dbmail/contrib/dbmailclient.php
usr/share/dbmail/contrib/exim-dbmail-configure
usr/share/dbmail/contrib/expiretrash.sh
usr/share/dbmail/contrib/mailbox2dbmail/
usr/share/dbmail/contrib/mailbox2dbmail/README
usr/share/dbmail/contrib/mailbox2dbmail/mailbox2dbmail
usr/share/dbmail/contrib/mailbox2dbmail/mailbox2dbmail.1
usr/share/dbmail/contrib/mailbox2dbmail/mailbox2dbmail.yo
usr/share/dbmail/contrib/nginx/
usr/share/dbmail/contrib/nginx/auth.php
usr/share/dbmail/contrib/nginx/mail.conf
usr/share/dbmail/contrib/sql2sql/
usr/share/dbmail/contrib/sql2sql/check_physmessage.sh
usr/share/dbmail/contrib/sql2sql/mysql2pgsql.sh
usr/share/dbmail/contrib/sql2sql/mysql2pgsql.yml
usr/share/dbmail/contrib/sql2sql/pgsql2mysql.sh
usr/share/dbmail/contrib/startup-scripts/
usr/share/dbmail/contrib/startup-scripts/README
usr/share/dbmail/contrib/startup-scripts/gentoo/
usr/share/dbmail/contrib/startup-scripts/gentoo/conf.d-dbmail
usr/share/dbmail/contrib/startup-scripts/gentoo/init.d-dbmail
usr/share/dbmail/contrib/startup-scripts/redhat/
usr/share/dbmail/contrib/startup-scripts/redhat/dbmail-imapd.init
usr/share/dbmail/contrib/startup-scripts/redhat/dbmail-lmtpd.init
usr/share/dbmail/contrib/startup-scripts/redhat/dbmail-pop3d.init
usr/share/dbmail/contrib/testdbmailclient.php
usr/share/dbmail/dbmail.schema
usr/share/dbmail/mysql/
usr/share/dbmail/mysql/2_1_7-2_2_0.mysql
usr/share/dbmail/mysql/2_2-3_0.mysql
usr/share/dbmail/mysql/2_2_5-2_3_0.mysql
usr/share/dbmail/mysql/2_3_0-2_3_2.mysql
usr/share/dbmail/mysql/2_3_4-2_3_5.mysql
usr/share/dbmail/mysql/2_3_5-2_3_6.mysql
usr/share/dbmail/mysql/2_3_6-3_0_0.mysql
usr/share/dbmail/mysql/3_0_2-3_1_0.mysql
usr/share/dbmail/mysql/create_tables.mysql
usr/share/dbmail/mysql/fix_foreign_keys.mysql
usr/share/dbmail/mysql/migrate_from_1.x_to_2.0_innodb.mysql
usr/share/dbmail/mysql/migrate_from_2.0_to_2.2.mysql
usr/share/dbmail/mysql/update_headervalue_01.mysql
usr/share/dbmail/mysql/upgrades/
usr/share/dbmail/mysql/upgrades/32001.mysql
usr/share/dbmail/mysql/upgrades/32002.mysql
usr/share/dbmail/mysql/upgrades/32003.mysql
usr/share/dbmail/mysql/upgrades/32004.mysql
usr/share/dbmail/mysql/upgrades/32005.mysql
usr/share/dbmail/mysql/upgrades/32006.mysql
usr/share/dbmail/mysql/upgrades/upgrade.mysql
usr/share/dbmail/oracle/
usr/share/dbmail/oracle/dbmail_grant.sql
usr/share/dbmail/oracle/dbmail_schema.sql
usr/share/dbmail/oracle/dbmail_ts.sql
usr/share/dbmail/postgresql/
usr/share/dbmail/postgresql/2_1_7-2_2_0.pgsql
usr/share/dbmail/postgresql/2_2-3_0.pgsql
usr/share/dbmail/postgresql/2_2_5-2_3_0.pgsql
usr/share/dbmail/postgresql/2_3_4-2_3_5.pgsql
usr/share/dbmail/postgresql/2_3_5-2_3_6.pgsql
usr/share/dbmail/postgresql/2_3_6-3_0_0.pgsql
usr/share/dbmail/postgresql/3_0_2-3_1_0.psql
usr/share/dbmail/postgresql/create_tables.pgsql
usr/share/dbmail/postgresql/migrate_from_1.x_to_2.0.pgsql
usr/share/dbmail/postgresql/migrate_from_2.0_to_2.2.pgsql
usr/share/dbmail/postgresql/upgrades/
usr/share/dbmail/postgresql/upgrades/32001.psql
usr/share/dbmail/postgresql/upgrades/32002.psql
usr/share/dbmail/postgresql/upgrades/32003.psql
usr/share/dbmail/postgresql/upgrades/32004.psql
usr/share/dbmail/postgresql/upgrades/32005.psql
usr/share/dbmail/postgresql/upgrades/32006.psql
usr/share/dbmail/postgresql/upgrades/upgrade.psql
usr/share/dbmail/sqlite/
usr/share/dbmail/sqlite/2_1_7-2_2_0.sqlite
usr/share/dbmail/sqlite/2_2-3_0.sqlite
usr/share/dbmail/sqlite/2_2_5-2_3_0.sqlite
usr/share/dbmail/sqlite/2_3_4-2_3_5.sqlite
usr/share/dbmail/sqlite/2_3_5-2_3_6.sqlite
usr/share/dbmail/sqlite/2_3_6-3_0_0.sqlite
usr/share/dbmail/sqlite/3_0_2-3_1_0.sqlite
usr/share/dbmail/sqlite/create_tables.sqlite
usr/share/dbmail/sqlite/trigger.tmpl.sql
usr/share/dbmail/sqlite/upgrades/
usr/share/dbmail/sqlite/upgrades/32001.sqlite
usr/share/dbmail/sqlite/upgrades/32002.sqlite
usr/share/dbmail/sqlite/upgrades/32003.sqlite
usr/share/dbmail/sqlite/upgrades/32004.sqlite
usr/share/dbmail/sqlite/upgrades/32005.sqlite
usr/share/dbmail/sqlite/upgrades/32006.sqlite
usr/share/dbmail/sqlite/upgrades/upgrade.sqlite
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dbmail-deliver.1.gz
usr/share/man/man1/dbmail.1.gz
usr/share/man/man5/
usr/share/man/man5/dbmail.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/dbmail-export.8.gz
usr/share/man/man8/dbmail-httpd.8.gz
usr/share/man/man8/dbmail-imapd.8.gz
usr/share/man/man8/dbmail-lmtpd.8.gz
usr/share/man/man8/dbmail-pop3d.8.gz
usr/share/man/man8/dbmail-sievecmd.8.gz
usr/share/man/man8/dbmail-sieved.8.gz
usr/share/man/man8/dbmail-timsieved.8.gz
usr/share/man/man8/dbmail-users.8.gz
usr/share/man/man8/dbmail-util.8.gz
var/
var/log/
var/log/dbmail/
