# extracted/prometheus-json-exporter-0.6.0-4/files
%FILES%
etc/
etc/prometheus/
etc/prometheus/json.yml
usr/
usr/bin/
usr/bin/prometheus-json-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-json-exporter.service
