# extracted/n2n-3.1.1-3/files
%FILES%
etc/
etc/n2n/
etc/n2n/edge.conf
etc/n2n/supernode.conf
usr/
usr/bin/
usr/bin/edge
usr/bin/n2n-benchmark
usr/bin/n2n-decode
usr/bin/n2n-keygen
usr/bin/supernode
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/edge.service
usr/lib/systemd/system/edge@.service
usr/lib/systemd/system/supernode.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/supernode.1.gz
usr/share/man/man7/
usr/share/man/man7/n2n.7.gz
usr/share/man/man8/
usr/share/man/man8/edge.8.gz
