# extracted/inn-2.7.2-3/files
%FILES%
etc/
etc/inn/
etc/inn/actsync.cfg
etc/inn/actsync.ign
etc/inn/buffindexed.conf
etc/inn/control.ctl
etc/inn/control.ctl.local
etc/inn/cycbuff.conf
etc/inn/distrib.pats
etc/inn/distributions
etc/inn/expire.ctl
etc/inn/incoming.conf
etc/inn/inn-radius.conf
etc/inn/inn-secrets.conf
etc/inn/inn.conf
etc/inn/innfeed.conf
etc/inn/innreport.conf
etc/inn/innshellvars.local
etc/inn/innshellvars.pl.local
etc/inn/innshellvars.tcl.local
etc/inn/innwatch.ctl
etc/inn/localgroups
etc/inn/moderators
etc/inn/motd.innd.sample
etc/inn/motd.nnrpd.sample
etc/inn/news2mail.cf
etc/inn/newsfeeds
etc/inn/nnrpd.track
etc/inn/nntpsend.ctl
etc/inn/nocem.ctl
etc/inn/ovdb.conf
etc/inn/ovsqlite.conf
etc/inn/passwd.nntp
etc/inn/readers.conf
etc/inn/send-uucp.cf
etc/inn/storage.conf
etc/inn/subscriptions
usr/
usr/bin/
usr/bin/actmerge
usr/bin/actsync
usr/bin/actsyncd
usr/bin/archive.inn
usr/bin/auth/
usr/bin/auth/passwd/
usr/bin/auth/passwd/auth_krb5
usr/bin/auth/passwd/ckpasswd
usr/bin/auth/passwd/radius
usr/bin/auth/resolv/
usr/bin/auth/resolv/domain
usr/bin/auth/resolv/ident
usr/bin/batcher
usr/bin/buffchan
usr/bin/buffindexed_d
usr/bin/cnfsheadconf
usr/bin/cnfsstat
usr/bin/control/
usr/bin/control/checkgroups.pl
usr/bin/control/ihave.pl
usr/bin/control/newgroup.pl
usr/bin/control/rmgroup.pl
usr/bin/control/sendme.pl
usr/bin/controlbatch
usr/bin/controlchan
usr/bin/convdate
usr/bin/ctlinnd
usr/bin/cvtbatch
usr/bin/delayer
usr/bin/docheckgroups
usr/bin/expire
usr/bin/expireover
usr/bin/expirerm
usr/bin/fastrm
usr/bin/filter/
usr/bin/filter/INN.py
usr/bin/filter/filter_innd.pl
usr/bin/filter/filter_innd.py
usr/bin/filter/filter_nnrpd.pl
usr/bin/filter/nnrpd.py
usr/bin/filter/nnrpd_access.pl
usr/bin/filter/nnrpd_access.py
usr/bin/filter/nnrpd_auth.pl
usr/bin/filter/nnrpd_auth.py
usr/bin/filter/nnrpd_dynamic.py
usr/bin/filter/startup_innd.pl
usr/bin/gencancel
usr/bin/getlist
usr/bin/grephistory
usr/bin/imapfeed
usr/bin/inews
usr/bin/innbind
usr/bin/inncheck
usr/bin/innconfval
usr/bin/innd
usr/bin/inndf
usr/bin/innfeed
usr/bin/innmail
usr/bin/innreport
usr/bin/innstat
usr/bin/innupgrade
usr/bin/innwatch
usr/bin/innxbatch
usr/bin/innxmit
usr/bin/mailpost
usr/bin/makedbz
usr/bin/makehistory
usr/bin/mod-active
usr/bin/news.daily
usr/bin/news2mail
usr/bin/ninpaths
usr/bin/nnrpd
usr/bin/nntpget
usr/bin/nntpsend
usr/bin/ovdb_init
usr/bin/ovdb_monitor
usr/bin/ovdb_server
usr/bin/ovdb_stat
usr/bin/overchan
usr/bin/ovsqlite-server
usr/bin/ovsqlite-util
usr/bin/perl-nocem
usr/bin/pgpverify
usr/bin/procbatch
usr/bin/prunehistory
usr/bin/pullnews
usr/bin/rc.news
usr/bin/rnews
usr/bin/rnews.libexec/
usr/bin/rnews.libexec/bunbatch
usr/bin/rnews.libexec/c7unbatch
usr/bin/rnews.libexec/decode
usr/bin/rnews.libexec/encode
usr/bin/rnews.libexec/gunbatch
usr/bin/scanlogs
usr/bin/scanspool
usr/bin/send-ihave
usr/bin/send-uucp
usr/bin/sendinpaths
usr/bin/sendxbatches
usr/bin/shlock
usr/bin/shrinkfile
usr/bin/simpleftp
usr/bin/sm
usr/bin/tally.control
usr/bin/tdx-util
usr/bin/tinyleaf
usr/bin/writelog
usr/include/
usr/include/inn/
usr/include/inn/inn/
usr/include/inn/inn/buffer.h
usr/include/inn/inn/concat.h
usr/include/inn/inn/confparse.h
usr/include/inn/inn/dbz.h
usr/include/inn/inn/dispatch.h
usr/include/inn/inn/fdflag.h
usr/include/inn/inn/hashtab.h
usr/include/inn/inn/history.h
usr/include/inn/inn/innconf.h
usr/include/inn/inn/inndcomm.h
usr/include/inn/inn/libinn.h
usr/include/inn/inn/list.h
usr/include/inn/inn/macros.h
usr/include/inn/inn/md5.h
usr/include/inn/inn/messages.h
usr/include/inn/inn/mmap.h
usr/include/inn/inn/network-innbind.h
usr/include/inn/inn/network.h
usr/include/inn/inn/newsuser.h
usr/include/inn/inn/nntp.h
usr/include/inn/inn/options.h
usr/include/inn/inn/ov.h
usr/include/inn/inn/overview.h
usr/include/inn/inn/paths.h
usr/include/inn/inn/portable-getaddrinfo.h
usr/include/inn/inn/portable-getnameinfo.h
usr/include/inn/inn/portable-macros.h
usr/include/inn/inn/portable-socket.h
usr/include/inn/inn/portable-stdbool.h
usr/include/inn/inn/qio.h
usr/include/inn/inn/secrets.h
usr/include/inn/inn/sequence.h
usr/include/inn/inn/storage.h
usr/include/inn/inn/system.h
usr/include/inn/inn/timer.h
usr/include/inn/inn/tst.h
usr/include/inn/inn/utility.h
usr/include/inn/inn/vector.h
usr/include/inn/inn/version.h
usr/include/inn/inn/wire.h
usr/include/inn/inn/xmalloc.h
usr/include/inn/inn/xwrite.h
usr/lib/
usr/lib/innreport-display.conf
usr/lib/innreport_inn.pm
usr/lib/innshellvars
usr/lib/innshellvars.pl
usr/lib/innshellvars.tcl
usr/lib/libinn.so
usr/lib/libinn.so.2
usr/lib/libinn.so.2.0
usr/lib/libinn.so.2.0.0
usr/lib/libinn.so.9
usr/lib/libinn.so.9.0.1
usr/lib/libinnhist.so
usr/lib/libinnhist.so.2
usr/lib/libinnhist.so.2.0
usr/lib/libinnhist.so.2.0.0
usr/lib/libinnhist.so.3
usr/lib/libinnhist.so.3.0.8
usr/lib/libinnstorage.so
usr/lib/libinnstorage.so.2
usr/lib/libinnstorage.so.2.0
usr/lib/libinnstorage.so.2.0.0
usr/lib/libinnstorage.so.3
usr/lib/libinnstorage.so.3.1.2
usr/lib/perl/
usr/lib/perl/INN/
usr/lib/perl/INN/Config.pm
usr/lib/perl/INN/Utils/
usr/lib/perl/INN/Utils/Shlock.pm
usr/lib/perl/INN/ovsqlite_client.pm
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/innd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/inn.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/inn.conf
usr/share/
usr/share/doc/
usr/share/doc/inn/
usr/share/doc/inn/CONTRIBUTORS
usr/share/doc/inn/GPL
usr/share/doc/inn/HACKING
usr/share/doc/inn/INSTALL
usr/share/doc/inn/LICENSE
usr/share/doc/inn/NEWS
usr/share/doc/inn/README
usr/share/doc/inn/TODO
usr/share/doc/inn/config-design
usr/share/doc/inn/config-semantics
usr/share/doc/inn/config-syntax
usr/share/doc/inn/external-auth
usr/share/doc/inn/history
usr/share/doc/inn/history-innfeed
usr/share/doc/inn/hook-perl
usr/share/doc/inn/hook-python
usr/share/doc/inn/sample-control
usr/share/inn/
usr/share/inn/innreport.css
usr/share/licenses/
usr/share/licenses/inn/
usr/share/licenses/inn/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/convdate.1.gz
usr/share/man/man1/delayer.1.gz
usr/share/man/man1/fastrm.1.gz
usr/share/man/man1/gencancel.1.gz
usr/share/man/man1/getlist.1.gz
usr/share/man/man1/grephistory.1.gz
usr/share/man/man1/inews.1.gz
usr/share/man/man1/innconfval.1.gz
usr/share/man/man1/innmail.1.gz
usr/share/man/man1/nntpget.1.gz
usr/share/man/man1/pgpverify.1.gz
usr/share/man/man1/pullnews.1.gz
usr/share/man/man1/rnews.1.gz
usr/share/man/man1/shlock.1.gz
usr/share/man/man1/shrinkfile.1.gz
usr/share/man/man1/simpleftp.1.gz
usr/share/man/man1/sm.1.gz
usr/share/man/man3/
usr/share/man/man3/INN::Config.3pm.gz
usr/share/man/man3/INN::Utils::Shlock.3pm.gz
usr/share/man/man3/INN::ovsqlite_client.3pm.gz
usr/share/man/man3/libinn.3.gz
usr/share/man/man3/libinn_clientlib.3.gz
usr/share/man/man3/libinn_dbz.3.gz
usr/share/man/man3/libinn_inndcomm.3.gz
usr/share/man/man3/libinn_list.3.gz
usr/share/man/man3/libinn_qio.3.gz
usr/share/man/man3/libinn_tst.3.gz
usr/share/man/man3/libinn_uwildmat.3.gz
usr/share/man/man3/libinnauth.3.gz
usr/share/man/man3/libinnhist.3.gz
usr/share/man/man3/libinnstorage.3.gz
usr/share/man/man5/
usr/share/man/man5/active.5.gz
usr/share/man/man5/active.times.5.gz
usr/share/man/man5/buffindexed.conf.5.gz
usr/share/man/man5/control.ctl.5.gz
usr/share/man/man5/cycbuff.conf.5.gz
usr/share/man/man5/distrib.pats.5.gz
usr/share/man/man5/distributions.5.gz
usr/share/man/man5/expire.ctl.5.gz
usr/share/man/man5/history.5.gz
usr/share/man/man5/incoming.conf.5.gz
usr/share/man/man5/inn-radius.conf.5.gz
usr/share/man/man5/inn-secrets.conf.5.gz
usr/share/man/man5/inn.conf.5.gz
usr/share/man/man5/innfeed.conf.5.gz
usr/share/man/man5/innreport.conf.5.gz
usr/share/man/man5/innwatch.ctl.5.gz
usr/share/man/man5/localgroups.5.gz
usr/share/man/man5/moderators.5.gz
usr/share/man/man5/motd.innd.5.gz
usr/share/man/man5/motd.news.5.gz
usr/share/man/man5/motd.nnrpd.5.gz
usr/share/man/man5/newsfeeds.5.gz
usr/share/man/man5/newsgroups.5.gz
usr/share/man/man5/newslog.5.gz
usr/share/man/man5/nnrpd.track.5.gz
usr/share/man/man5/nntpsend.ctl.5.gz
usr/share/man/man5/nocem.ctl.5.gz
usr/share/man/man5/ovdb.5.gz
usr/share/man/man5/ovsqlite.5.gz
usr/share/man/man5/passwd.nntp.5.gz
usr/share/man/man5/readers.conf.5.gz
usr/share/man/man5/storage.conf.5.gz
usr/share/man/man5/subscriptions.5.gz
usr/share/man/man8/
usr/share/man/man8/actsync.8.gz
usr/share/man/man8/actsyncd.8.gz
usr/share/man/man8/archive.8.gz
usr/share/man/man8/auth_krb5.8.gz
usr/share/man/man8/batcher.8.gz
usr/share/man/man8/buffchan.8.gz
usr/share/man/man8/ckpasswd.8.gz
usr/share/man/man8/cnfsheadconf.8.gz
usr/share/man/man8/cnfsstat.8.gz
usr/share/man/man8/controlchan.8.gz
usr/share/man/man8/ctlinnd.8.gz
usr/share/man/man8/cvtbatch.8.gz
usr/share/man/man8/docheckgroups.8.gz
usr/share/man/man8/domain.8.gz
usr/share/man/man8/expire.8.gz
usr/share/man/man8/expireover.8.gz
usr/share/man/man8/expirerm.8.gz
usr/share/man/man8/ident.8.gz
usr/share/man/man8/imapfeed.8.gz
usr/share/man/man8/innbind.8.gz
usr/share/man/man8/inncheck.8.gz
usr/share/man/man8/innd.8.gz
usr/share/man/man8/inndf.8.gz
usr/share/man/man8/innfeed.8.gz
usr/share/man/man8/innreport.8.gz
usr/share/man/man8/innstat.8.gz
usr/share/man/man8/innupgrade.8.gz
usr/share/man/man8/innwatch.8.gz
usr/share/man/man8/innxbatch.8.gz
usr/share/man/man8/innxmit.8.gz
usr/share/man/man8/inpaths.8.gz
usr/share/man/man8/mailpost.8.gz
usr/share/man/man8/makedbz.8.gz
usr/share/man/man8/makehistory.8.gz
usr/share/man/man8/mod-active.8.gz
usr/share/man/man8/news.daily.8.gz
usr/share/man/man8/news2mail.8.gz
usr/share/man/man8/ninpaths.8.gz
usr/share/man/man8/nnrpd.8.gz
usr/share/man/man8/nntpsend.8.gz
usr/share/man/man8/ovdb_init.8.gz
usr/share/man/man8/ovdb_monitor.8.gz
usr/share/man/man8/ovdb_server.8.gz
usr/share/man/man8/ovdb_stat.8.gz
usr/share/man/man8/overchan.8.gz
usr/share/man/man8/ovsqlite-server.8.gz
usr/share/man/man8/ovsqlite-util.8.gz
usr/share/man/man8/perl-nocem.8.gz
usr/share/man/man8/procbatch.8.gz
usr/share/man/man8/prunehistory.8.gz
usr/share/man/man8/radius.8.gz
usr/share/man/man8/rc.news.8.gz
usr/share/man/man8/scanlogs.8.gz
usr/share/man/man8/scanspool.8.gz
usr/share/man/man8/send-ihave.8.gz
usr/share/man/man8/send-uucp.8.gz
usr/share/man/man8/sendinpaths.8.gz
usr/share/man/man8/sendxbatches.8.gz
usr/share/man/man8/tally.control.8.gz
usr/share/man/man8/tdx-util.8.gz
usr/share/man/man8/tinyleaf.8.gz
usr/share/man/man8/writelog.8.gz
var/
var/db/
var/db/inn/
var/db/inn/active
var/db/inn/active.times
var/db/inn/history
var/db/inn/history.dir
var/db/inn/history.hash
var/db/inn/history.index
var/db/inn/newsgroups
var/log/
var/log/inn/
var/log/inn/OLD/
var/spool/
var/spool/inn/
var/spool/inn/tmp/
var/spool/news/
var/spool/news/archive/
var/spool/news/articles/
var/spool/news/incoming/
var/spool/news/incoming/bad/
var/spool/news/innfeed/
var/spool/news/outgoing/
var/spool/news/overview/
