# extracted/qemu-guest-agent-9.1.2-1/files
%FILES%
etc/
etc/qemu/
etc/qemu/fsfreeze-hook
etc/qemu/fsfreeze-hook.d/
etc/qemu/qemu-ga.conf
usr/
usr/bin/
usr/bin/qemu-ga
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/qemu-guest-agent.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-qemu-guest-agent.rules
usr/share/
usr/share/licenses/
usr/share/licenses/qemu-guest-agent/
usr/share/licenses/qemu-guest-agent/BSD-2-Clause.txt
usr/share/licenses/qemu-guest-agent/LICENSE
usr/share/licenses/qemu-guest-agent/audio.LICENSE.ISC.txt
usr/share/licenses/qemu-guest-agent/edk2-licenses.txt
usr/share/licenses/qemu-guest-agent/fuse.LICENSE.BSD-2-Clause.txt
usr/share/licenses/qemu-guest-agent/tcg-arm.LICENSE.MIT.txt
usr/share/licenses/qemu-guest-agent/tcg.LICENSE.MIT.txt
usr/share/licenses/qemu-guest-agent/tci.LICENSE.MIT.txt
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/qemu-ga.8.gz
