# extracted/scx-scheds-1.0.7-1/files
%FILES%
etc/
etc/default/
etc/default/scx
usr/
usr/bin/
usr/bin/scx_bpfland
usr/bin/scx_central
usr/bin/scx_flash
usr/bin/scx_flatcg
usr/bin/scx_lavd
usr/bin/scx_layered
usr/bin/scx_loader
usr/bin/scx_nest
usr/bin/scx_pair
usr/bin/scx_qmap
usr/bin/scx_rlfifo
usr/bin/scx_rustland
usr/bin/scx_rusty
usr/bin/scx_sdt
usr/bin/scx_simple
usr/bin/scx_userland
usr/include/
usr/include/scx/
usr/include/scx/bpf_arena_common.h
usr/include/scx/common.bpf.h
usr/include/scx/common.h
usr/include/scx/compat.bpf.h
usr/include/scx/compat.h
usr/include/scx/enums.autogen.bpf.h
usr/include/scx/enums.autogen.h
usr/include/scx/enums.bpf.h
usr/include/scx/enums.h
usr/include/scx/namespace.bpf.h
usr/include/scx/namespace_impl.bpf.h
usr/include/scx/ravg.bpf.h
usr/include/scx/ravg_impl.bpf.h
usr/include/scx/user_exit_info.h
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/scx.service
usr/lib/systemd/system/scx_loader.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.scx.Loader.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.scx.Loader.conf
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/90-scx-scheds-upgrade.hook
usr/share/libalpm/scripts/
usr/share/libalpm/scripts/scx-systemd-restart
