# extracted/recoll-1.41.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/recoll
usr/bin/recollindex
usr/bin/recollq
usr/include/
usr/include/recoll/
usr/include/recoll/Filter.h
usr/include/recoll/docseq.h
usr/include/recoll/docseqdb.h
usr/include/recoll/hldata.h
usr/include/recoll/idxstatus.h
usr/include/recoll/internfile.h
usr/include/recoll/log.h
usr/include/recoll/mimehandler.h
usr/include/recoll/pathut.h
usr/include/recoll/plaintorich.h
usr/include/recoll/qresultstore.h
usr/include/recoll/rclconfig.h
usr/include/recoll/rcldb.h
usr/include/recoll/rcldoc.h
usr/include/recoll/rclinit.h
usr/include/recoll/rclquery.h
usr/include/recoll/rclutil.h
usr/include/recoll/readfile.h
usr/include/recoll/reslistpager.h
usr/include/recoll/searchdata.h
usr/include/recoll/smallut.h
usr/include/recoll/wasatorcl.h
usr/lib/
usr/lib/librecoll.so
usr/lib/librecoll.so.39
usr/lib/librecoll.so.39.2.0
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/recoll/
usr/lib/python3.12/site-packages/recoll/__init__.py
usr/lib/python3.12/site-packages/recoll/__pycache__/
usr/lib/python3.12/site-packages/recoll/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/conftree.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/fsudi.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/qresultstore.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/rclconfig.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/rclextract.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/__pycache__/recoll.cpython-312.pyc
usr/lib/python3.12/site-packages/recoll/_recoll.cpython-312-x86_64-linux-gnu.so
usr/lib/python3.12/site-packages/recoll/conftree.py
usr/lib/python3.12/site-packages/recoll/fsudi.py
usr/lib/python3.12/site-packages/recoll/qresultstore.py
usr/lib/python3.12/site-packages/recoll/rclconfig.py
usr/lib/python3.12/site-packages/recoll/rclextract.py
usr/lib/python3.12/site-packages/recoll/recoll.py
usr/lib/python3.12/site-packages/recollaspell.cpython-312-x86_64-linux-gnu.so
usr/lib/python3.12/site-packages/recollchm/
usr/lib/python3.12/site-packages/recollchm/__init__.py
usr/lib/python3.12/site-packages/recollchm/__pycache__/
usr/lib/python3.12/site-packages/recollchm/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/recollchm/__pycache__/chm.cpython-312.pyc
usr/lib/python3.12/site-packages/recollchm/__pycache__/chmlib.cpython-312.pyc
usr/lib/python3.12/site-packages/recollchm/_chmlib.cpython-312-x86_64-linux-gnu.so
usr/lib/python3.12/site-packages/recollchm/chm.py
usr/lib/python3.12/site-packages/recollchm/chmlib.py
usr/lib/qt6/
usr/lib/qt6/plugins/
usr/lib/qt6/plugins/kf6/
usr/lib/qt6/plugins/kf6/kio/
usr/lib/qt6/plugins/kf6/kio/kio_recoll.so
usr/lib/qt6/plugins/kf6/krunner/
usr/lib/qt6/plugins/kf6/krunner/krunner_recoll.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/recollindex@.service
usr/lib/systemd/user/
usr/lib/systemd/user/recollindex.service
usr/share/
usr/share/applications/
usr/share/applications/recoll-searchgui.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/apps/
usr/share/icons/hicolor/48x48/apps/recoll.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/recoll.svg
usr/share/kio_recoll/
usr/share/kio_recoll/help.html
usr/share/kio_recoll/welcome.html
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/recoll.1.gz
usr/share/man/man1/recollindex.1.gz
usr/share/man/man1/recollq.1.gz
usr/share/man/man1/xadump.1.gz
usr/share/man/man5/
usr/share/man/man5/recoll.conf.5.gz
usr/share/metainfo/
usr/share/metainfo/org.recoll.recoll.appdata.xml
usr/share/pixmaps/
usr/share/pixmaps/recoll.png
usr/share/recoll/
usr/share/recoll/doc/
usr/share/recoll/doc/docbook-xsl.css
usr/share/recoll/doc/usermanual.html
usr/share/recoll/examples/
usr/share/recoll/examples/backends
usr/share/recoll/examples/fields
usr/share/recoll/examples/fragment-buttons.xml
usr/share/recoll/examples/mimeconf
usr/share/recoll/examples/mimemap
usr/share/recoll/examples/mimeview
usr/share/recoll/examples/rclmon.sh
usr/share/recoll/examples/recoll-common.css
usr/share/recoll/examples/recoll-common.qss
usr/share/recoll/examples/recoll-dark.css
usr/share/recoll/examples/recoll-dark.qss
usr/share/recoll/examples/recoll.conf
usr/share/recoll/examples/recoll.qss
usr/share/recoll/examples/recollindex.desktop
usr/share/recoll/examples/recollindex.service
usr/share/recoll/examples/recollindex@.service
usr/share/recoll/filters/
usr/share/recoll/filters/abiword.xsl
usr/share/recoll/filters/archivextract.py
usr/share/recoll/filters/cmdtalk.py
usr/share/recoll/filters/cnsplitter.py
usr/share/recoll/filters/conftree.py
usr/share/recoll/filters/fb2.xsl
usr/share/recoll/filters/gnumeric.xsl
usr/share/recoll/filters/kosplitter.py
usr/share/recoll/filters/msodump.zip
usr/share/recoll/filters/okular-note.xsl
usr/share/recoll/filters/opendoc-body.xsl
usr/share/recoll/filters/opendoc-flat.xsl
usr/share/recoll/filters/opendoc-meta.xsl
usr/share/recoll/filters/openxml-meta.xsl
usr/share/recoll/filters/openxml-word-body.xsl
usr/share/recoll/filters/openxml-xls-body.xsl
usr/share/recoll/filters/ppt-dump.py
usr/share/recoll/filters/rcl7z.py
usr/share/recoll/filters/rclaptosidman
usr/share/recoll/filters/rclaspell-sugg.py
usr/share/recoll/filters/rclaudio.py
usr/share/recoll/filters/rclbasehandler.py
usr/share/recoll/filters/rclbibtex.sh
usr/share/recoll/filters/rclcheckneedretry.sh
usr/share/recoll/filters/rclchm.py
usr/share/recoll/filters/rclconfig.py
usr/share/recoll/filters/rcldia.py
usr/share/recoll/filters/rcldjvu.py
usr/share/recoll/filters/rcldoc.py
usr/share/recoll/filters/rcldvi
usr/share/recoll/filters/rclepub.py
usr/share/recoll/filters/rclepub1.py
usr/share/recoll/filters/rclexec1.py
usr/share/recoll/filters/rclexecm.py
usr/share/recoll/filters/rclfb2.py
usr/share/recoll/filters/rclgaim
usr/share/recoll/filters/rclgenxslt.py
usr/share/recoll/filters/rclhwp.py
usr/share/recoll/filters/rclics.py
usr/share/recoll/filters/rclimg
usr/share/recoll/filters/rclimg.py
usr/share/recoll/filters/rclinfo.py
usr/share/recoll/filters/rclipynb.py
usr/share/recoll/filters/rcljoplin.py
usr/share/recoll/filters/rclkar.py
usr/share/recoll/filters/rclkwd
usr/share/recoll/filters/rcllatinclass.py
usr/share/recoll/filters/rcllatinstops.zip
usr/share/recoll/filters/rcllink.py
usr/share/recoll/filters/rcllyx
usr/share/recoll/filters/rclman
usr/share/recoll/filters/rclmd.py
usr/share/recoll/filters/rclmidi.py
usr/share/recoll/filters/rclnamefilter.py
usr/share/recoll/filters/rclnull.py
usr/share/recoll/filters/rclocr.py
usr/share/recoll/filters/rclocrabbyy.py
usr/share/recoll/filters/rclocrcache.py
usr/share/recoll/filters/rclocrtesseract.py
usr/share/recoll/filters/rclopxml.py
usr/share/recoll/filters/rclorgmode.py
usr/share/recoll/filters/rclpdf.py
usr/share/recoll/filters/rclppt.py
usr/share/recoll/filters/rclps
usr/share/recoll/filters/rclpst.py
usr/share/recoll/filters/rclpurple
usr/share/recoll/filters/rclpython.py
usr/share/recoll/filters/rclrar.py
usr/share/recoll/filters/rclrtf.py
usr/share/recoll/filters/rclscribus
usr/share/recoll/filters/rclsfz.py
usr/share/recoll/filters/rclshowinfo
usr/share/recoll/filters/rcltar.py
usr/share/recoll/filters/rcltex
usr/share/recoll/filters/rcltext.py
usr/share/recoll/filters/rcltxtlines.py
usr/share/recoll/filters/rcluncomp
usr/share/recoll/filters/rcluncomp.py
usr/share/recoll/filters/rclwar.py
usr/share/recoll/filters/rclxls.py
usr/share/recoll/filters/rclxml.py
usr/share/recoll/filters/rclxmp.py
usr/share/recoll/filters/rclxslt.py
usr/share/recoll/filters/rclzip.py
usr/share/recoll/filters/recoll-we-move-files.py
usr/share/recoll/filters/recollepub.zip
usr/share/recoll/filters/svg.xsl
usr/share/recoll/filters/thunderbird-open-message.sh
usr/share/recoll/filters/xls-dump.py
usr/share/recoll/filters/xlsxmltocsv.py
usr/share/recoll/filters/xml.xsl
usr/share/recoll/images/
usr/share/recoll/images/aptosid-book.png
usr/share/recoll/images/aptosid-manual.png
usr/share/recoll/images/archive.png
usr/share/recoll/images/book.png
usr/share/recoll/images/bookchap.png
usr/share/recoll/images/document.png
usr/share/recoll/images/drawing.png
usr/share/recoll/images/emblem-symbolic-link.png
usr/share/recoll/images/folder.png
usr/share/recoll/images/html.png
usr/share/recoll/images/image.png
usr/share/recoll/images/message.png
usr/share/recoll/images/mozilla_doc.png
usr/share/recoll/images/pdf.png
usr/share/recoll/images/pidgin.png
usr/share/recoll/images/postscript.png
usr/share/recoll/images/presentation.png
usr/share/recoll/images/sidux-book.png
usr/share/recoll/images/soffice.png
usr/share/recoll/images/source.png
usr/share/recoll/images/sownd.png
usr/share/recoll/images/spreadsheet.png
usr/share/recoll/images/text-x-python.png
usr/share/recoll/images/txt.png
usr/share/recoll/images/video.png
usr/share/recoll/images/wordprocessing.png
usr/share/recoll/translations/
usr/share/recoll/translations/recoll_cs.qm
usr/share/recoll/translations/recoll_da.qm
usr/share/recoll/translations/recoll_de.qm
usr/share/recoll/translations/recoll_el.qm
usr/share/recoll/translations/recoll_es.qm
usr/share/recoll/translations/recoll_fr.qm
usr/share/recoll/translations/recoll_hu.qm
usr/share/recoll/translations/recoll_it.qm
usr/share/recoll/translations/recoll_ja.qm
usr/share/recoll/translations/recoll_ko.qm
usr/share/recoll/translations/recoll_lt.qm
usr/share/recoll/translations/recoll_nl.qm
usr/share/recoll/translations/recoll_pl.qm
usr/share/recoll/translations/recoll_ru.qm
usr/share/recoll/translations/recoll_sv.qm
usr/share/recoll/translations/recoll_tr.qm
usr/share/recoll/translations/recoll_uk.qm
usr/share/recoll/translations/recoll_xx.qm
usr/share/recoll/translations/recoll_zh.qm
usr/share/recoll/translations/recoll_zh_CN.qm
