# extracted/sysprof-47.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/sysprof
usr/bin/sysprof-agent
usr/bin/sysprof-cli
usr/include/
usr/include/sysprof-6/
usr/include/sysprof-6/sysprof-battery-charge.h
usr/include/sysprof-6/sysprof-bundled-symbolizer.h
usr/include/sysprof-6/sysprof-callgraph-frame.h
usr/include/sysprof-6/sysprof-callgraph-symbol.h
usr/include/sysprof-6/sysprof-callgraph.h
usr/include/sysprof-6/sysprof-category-summary.h
usr/include/sysprof-6/sysprof-cpu-info.h
usr/include/sysprof-6/sysprof-cpu-usage.h
usr/include/sysprof-6/sysprof-dbus-monitor.h
usr/include/sysprof-6/sysprof-diagnostic.h
usr/include/sysprof-6/sysprof-disk-usage.h
usr/include/sysprof-6/sysprof-document-allocation.h
usr/include/sysprof-6/sysprof-document-counter-value.h
usr/include/sysprof-6/sysprof-document-counter.h
usr/include/sysprof-6/sysprof-document-ctrdef.h
usr/include/sysprof-6/sysprof-document-ctrset.h
usr/include/sysprof-6/sysprof-document-dbus-message.h
usr/include/sysprof-6/sysprof-document-exit.h
usr/include/sysprof-6/sysprof-document-file-chunk.h
usr/include/sysprof-6/sysprof-document-file.h
usr/include/sysprof-6/sysprof-document-fork.h
usr/include/sysprof-6/sysprof-document-frame.h
usr/include/sysprof-6/sysprof-document-jitmap.h
usr/include/sysprof-6/sysprof-document-loader.h
usr/include/sysprof-6/sysprof-document-log.h
usr/include/sysprof-6/sysprof-document-mark.h
usr/include/sysprof-6/sysprof-document-metadata.h
usr/include/sysprof-6/sysprof-document-mmap.h
usr/include/sysprof-6/sysprof-document-overlay.h
usr/include/sysprof-6/sysprof-document-process.h
usr/include/sysprof-6/sysprof-document-sample.h
usr/include/sysprof-6/sysprof-document-traceable.h
usr/include/sysprof-6/sysprof-document.h
usr/include/sysprof-6/sysprof-elf-symbolizer.h
usr/include/sysprof-6/sysprof-energy-usage.h
usr/include/sysprof-6/sysprof-enums.h
usr/include/sysprof-6/sysprof-instrument.h
usr/include/sysprof-6/sysprof-jitmap-symbolizer.h
usr/include/sysprof-6/sysprof-kallsyms-symbolizer.h
usr/include/sysprof-6/sysprof-malloc-tracing.h
usr/include/sysprof-6/sysprof-mark-catalog.h
usr/include/sysprof-6/sysprof-memory-usage.h
usr/include/sysprof-6/sysprof-mount.h
usr/include/sysprof-6/sysprof-multi-symbolizer.h
usr/include/sysprof-6/sysprof-network-usage.h
usr/include/sysprof-6/sysprof-no-symbolizer.h
usr/include/sysprof-6/sysprof-power-profile.h
usr/include/sysprof-6/sysprof-profiler.h
usr/include/sysprof-6/sysprof-proxied-instrument.h
usr/include/sysprof-6/sysprof-recording.h
usr/include/sysprof-6/sysprof-sampler.h
usr/include/sysprof-6/sysprof-scheduler-details.h
usr/include/sysprof-6/sysprof-spawnable.h
usr/include/sysprof-6/sysprof-subprocess-output.h
usr/include/sysprof-6/sysprof-symbol.h
usr/include/sysprof-6/sysprof-symbolizer.h
usr/include/sysprof-6/sysprof-symbols-bundle.h
usr/include/sysprof-6/sysprof-system-logs.h
usr/include/sysprof-6/sysprof-thread-info.h
usr/include/sysprof-6/sysprof-time-span.h
usr/include/sysprof-6/sysprof-tracefd-consumer.h
usr/include/sysprof-6/sysprof-tracer.h
usr/include/sysprof-6/sysprof.h
usr/lib/
usr/lib/libsysprof-6.so
usr/lib/libsysprof-6.so.6
usr/lib/libsysprof-6.so.6.0.0
usr/lib/libsysprof-memory-6.so
usr/lib/libsysprof-speedtrack-6.so
usr/lib/libsysprof-tracer-6.so
usr/lib/pkgconfig/
usr/lib/pkgconfig/sysprof-6.pc
usr/lib/sysprofd
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sysprof3.service
usr/share/
usr/share/applications/
usr/share/applications/org.gnome.Sysprof.desktop
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.gnome.Sysprof.Agent.xml
usr/share/dbus-1/interfaces/org.gnome.Sysprof3.Profiler.xml
usr/share/dbus-1/interfaces/org.gnome.Sysprof3.Service.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.gnome.Sysprof3.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.gnome.Sysprof3.conf
usr/share/help/
usr/share/help/C/
usr/share/help/C/sysprof/
usr/share/help/C/sysprof/faq.page
usr/share/help/C/sysprof/figures/
usr/share/help/C/sysprof/figures/Sysprof.svg
usr/share/help/C/sysprof/index.page
usr/share/help/C/sysprof/introduction.page
usr/share/help/C/sysprof/legal.xml
usr/share/help/C/sysprof/profiling.page
usr/share/help/ca/
usr/share/help/ca/sysprof/
usr/share/help/ca/sysprof/faq.page
usr/share/help/ca/sysprof/figures/
usr/share/help/ca/sysprof/figures/Sysprof.svg
usr/share/help/ca/sysprof/index.page
usr/share/help/ca/sysprof/introduction.page
usr/share/help/ca/sysprof/legal.xml
usr/share/help/ca/sysprof/profiling.page
usr/share/help/cs/
usr/share/help/cs/sysprof/
usr/share/help/cs/sysprof/faq.page
usr/share/help/cs/sysprof/figures/
usr/share/help/cs/sysprof/figures/Sysprof.svg
usr/share/help/cs/sysprof/index.page
usr/share/help/cs/sysprof/introduction.page
usr/share/help/cs/sysprof/legal.xml
usr/share/help/cs/sysprof/profiling.page
usr/share/help/da/
usr/share/help/da/sysprof/
usr/share/help/da/sysprof/faq.page
usr/share/help/da/sysprof/figures/
usr/share/help/da/sysprof/figures/Sysprof.svg
usr/share/help/da/sysprof/index.page
usr/share/help/da/sysprof/introduction.page
usr/share/help/da/sysprof/legal.xml
usr/share/help/da/sysprof/profiling.page
usr/share/help/de/
usr/share/help/de/sysprof/
usr/share/help/de/sysprof/faq.page
usr/share/help/de/sysprof/figures/
usr/share/help/de/sysprof/figures/Sysprof.svg
usr/share/help/de/sysprof/index.page
usr/share/help/de/sysprof/introduction.page
usr/share/help/de/sysprof/legal.xml
usr/share/help/de/sysprof/profiling.page
usr/share/help/el/
usr/share/help/el/sysprof/
usr/share/help/el/sysprof/faq.page
usr/share/help/el/sysprof/figures/
usr/share/help/el/sysprof/figures/Sysprof.svg
usr/share/help/el/sysprof/index.page
usr/share/help/el/sysprof/introduction.page
usr/share/help/el/sysprof/legal.xml
usr/share/help/el/sysprof/profiling.page
usr/share/help/es/
usr/share/help/es/sysprof/
usr/share/help/es/sysprof/faq.page
usr/share/help/es/sysprof/figures/
usr/share/help/es/sysprof/figures/Sysprof.svg
usr/share/help/es/sysprof/index.page
usr/share/help/es/sysprof/introduction.page
usr/share/help/es/sysprof/legal.xml
usr/share/help/es/sysprof/profiling.page
usr/share/help/eu/
usr/share/help/eu/sysprof/
usr/share/help/eu/sysprof/faq.page
usr/share/help/eu/sysprof/figures/
usr/share/help/eu/sysprof/figures/Sysprof.svg
usr/share/help/eu/sysprof/index.page
usr/share/help/eu/sysprof/introduction.page
usr/share/help/eu/sysprof/legal.xml
usr/share/help/eu/sysprof/profiling.page
usr/share/help/fr/
usr/share/help/fr/sysprof/
usr/share/help/fr/sysprof/faq.page
usr/share/help/fr/sysprof/figures/
usr/share/help/fr/sysprof/figures/Sysprof.svg
usr/share/help/fr/sysprof/index.page
usr/share/help/fr/sysprof/introduction.page
usr/share/help/fr/sysprof/legal.xml
usr/share/help/fr/sysprof/profiling.page
usr/share/help/hu/
usr/share/help/hu/sysprof/
usr/share/help/hu/sysprof/faq.page
usr/share/help/hu/sysprof/figures/
usr/share/help/hu/sysprof/figures/Sysprof.svg
usr/share/help/hu/sysprof/index.page
usr/share/help/hu/sysprof/introduction.page
usr/share/help/hu/sysprof/legal.xml
usr/share/help/hu/sysprof/profiling.page
usr/share/help/id/
usr/share/help/id/sysprof/
usr/share/help/id/sysprof/faq.page
usr/share/help/id/sysprof/figures/
usr/share/help/id/sysprof/figures/Sysprof.svg
usr/share/help/id/sysprof/index.page
usr/share/help/id/sysprof/introduction.page
usr/share/help/id/sysprof/legal.xml
usr/share/help/id/sysprof/profiling.page
usr/share/help/ko/
usr/share/help/ko/sysprof/
usr/share/help/ko/sysprof/faq.page
usr/share/help/ko/sysprof/figures/
usr/share/help/ko/sysprof/figures/Sysprof.svg
usr/share/help/ko/sysprof/index.page
usr/share/help/ko/sysprof/introduction.page
usr/share/help/ko/sysprof/legal.xml
usr/share/help/ko/sysprof/profiling.page
usr/share/help/pl/
usr/share/help/pl/sysprof/
usr/share/help/pl/sysprof/faq.page
usr/share/help/pl/sysprof/figures/
usr/share/help/pl/sysprof/figures/Sysprof.svg
usr/share/help/pl/sysprof/index.page
usr/share/help/pl/sysprof/introduction.page
usr/share/help/pl/sysprof/legal.xml
usr/share/help/pl/sysprof/profiling.page
usr/share/help/pt_BR/
usr/share/help/pt_BR/sysprof/
usr/share/help/pt_BR/sysprof/faq.page
usr/share/help/pt_BR/sysprof/figures/
usr/share/help/pt_BR/sysprof/figures/Sysprof.svg
usr/share/help/pt_BR/sysprof/index.page
usr/share/help/pt_BR/sysprof/introduction.page
usr/share/help/pt_BR/sysprof/legal.xml
usr/share/help/pt_BR/sysprof/profiling.page
usr/share/help/ru/
usr/share/help/ru/sysprof/
usr/share/help/ru/sysprof/faq.page
usr/share/help/ru/sysprof/figures/
usr/share/help/ru/sysprof/figures/Sysprof.svg
usr/share/help/ru/sysprof/index.page
usr/share/help/ru/sysprof/introduction.page
usr/share/help/ru/sysprof/legal.xml
usr/share/help/ru/sysprof/profiling.page
usr/share/help/sv/
usr/share/help/sv/sysprof/
usr/share/help/sv/sysprof/faq.page
usr/share/help/sv/sysprof/figures/
usr/share/help/sv/sysprof/figures/Sysprof.svg
usr/share/help/sv/sysprof/index.page
usr/share/help/sv/sysprof/introduction.page
usr/share/help/sv/sysprof/legal.xml
usr/share/help/sv/sysprof/profiling.page
usr/share/help/uk/
usr/share/help/uk/sysprof/
usr/share/help/uk/sysprof/faq.page
usr/share/help/uk/sysprof/figures/
usr/share/help/uk/sysprof/figures/Sysprof.svg
usr/share/help/uk/sysprof/index.page
usr/share/help/uk/sysprof/introduction.page
usr/share/help/uk/sysprof/legal.xml
usr/share/help/uk/sysprof/profiling.page
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/actions/
usr/share/icons/hicolor/scalable/actions/sysprof-allocations.svg
usr/share/icons/hicolor/scalable/actions/sysprof-battery.svg
usr/share/icons/hicolor/scalable/actions/sysprof-calgraph.svg
usr/share/icons/hicolor/scalable/actions/sysprof-cli.svg
usr/share/icons/hicolor/scalable/actions/sysprof-cpu.svg
usr/share/icons/hicolor/scalable/actions/sysprof-disk.svg
usr/share/icons/hicolor/scalable/actions/sysprof-gtk.svg
usr/share/icons/hicolor/scalable/actions/sysprof-library.svg
usr/share/icons/hicolor/scalable/actions/sysprof-memory.svg
usr/share/icons/hicolor/scalable/actions/sysprof-networking.svg
usr/share/icons/hicolor/scalable/actions/sysprof-rapl.svg
usr/share/icons/hicolor/scalable/actions/sysprof-trace-app.svg
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/org.gnome.Sysprof.svg
usr/share/icons/hicolor/symbolic/
usr/share/icons/hicolor/symbolic/apps/
usr/share/icons/hicolor/symbolic/apps/org.gnome.Sysprof-symbolic.svg
usr/share/locale/
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/sysprof.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/sysprof.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/sysprof.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/sysprof.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/sysprof.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/sysprof.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/sysprof.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/sysprof.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/sysprof.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/sysprof.mo
usr/share/locale/fa/
usr/share/locale/fa/LC_MESSAGES/
usr/share/locale/fa/LC_MESSAGES/sysprof.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/sysprof.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/sysprof.mo
usr/share/locale/fur/
usr/share/locale/fur/LC_MESSAGES/
usr/share/locale/fur/LC_MESSAGES/sysprof.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/sysprof.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/sysprof.mo
usr/share/locale/hi/
usr/share/locale/hi/LC_MESSAGES/
usr/share/locale/hi/LC_MESSAGES/sysprof.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/sysprof.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/sysprof.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/sysprof.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/sysprof.mo
usr/share/locale/kab/
usr/share/locale/kab/LC_MESSAGES/
usr/share/locale/kab/LC_MESSAGES/sysprof.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/sysprof.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/sysprof.mo
usr/share/locale/ne/
usr/share/locale/ne/LC_MESSAGES/
usr/share/locale/ne/LC_MESSAGES/sysprof.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/sysprof.mo
usr/share/locale/oc/
usr/share/locale/oc/LC_MESSAGES/
usr/share/locale/oc/LC_MESSAGES/sysprof.mo
usr/share/locale/pa/
usr/share/locale/pa/LC_MESSAGES/
usr/share/locale/pa/LC_MESSAGES/sysprof.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/sysprof.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/sysprof.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/sysprof.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/sysprof.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/sysprof.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/sysprof.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/sysprof.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/sysprof.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/sysprof.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/sysprof.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/sysprof.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/sysprof.mo
usr/share/metainfo/
usr/share/metainfo/org.gnome.Sysprof.appdata.xml
usr/share/mime/
usr/share/mime/packages/
usr/share/mime/packages/sysprof-mime.xml
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.gnome.sysprof3.policy
