# extracted/tuntox-0.0.10.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/tokssh
usr/bin/tuntox
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tuntox.service
