# extracted/python-location-0.9.17-4/files
%FILES%
usr/
usr/bin/
usr/bin/location
usr/bin/location-importer
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/_location.so
usr/lib/python3.12/site-packages/location/
usr/lib/python3.12/site-packages/location/__init__.py
usr/lib/python3.12/site-packages/location/__pycache__/
usr/lib/python3.12/site-packages/location/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/database.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/database.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/downloader.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/downloader.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/export.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/export.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/i18n.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/i18n.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/importer.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/importer.cpython-312.pyc
usr/lib/python3.12/site-packages/location/__pycache__/logger.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/location/__pycache__/logger.cpython-312.pyc
usr/lib/python3.12/site-packages/location/database.py
usr/lib/python3.12/site-packages/location/downloader.py
usr/lib/python3.12/site-packages/location/export.py
usr/lib/python3.12/site-packages/location/i18n.py
usr/lib/python3.12/site-packages/location/importer.py
usr/lib/python3.12/site-packages/location/logger.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/location-update.service
usr/lib/systemd/system/location-update.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/location
usr/share/licenses/
usr/share/licenses/python-location/
usr/share/licenses/python-location/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/location.1.gz
var/
var/lib/
var/lib/location/
var/lib/location/signing-key.pem
