# extracted/quassel-core-0.14.0-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/quassel
usr/
usr/bin/
usr/bin/quasselcore
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/quassel.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/quassel.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/quassel.conf
