# extracted/openntpd-6.8p1-8/files
%FILES%
etc/
etc/ntpd.conf
usr/
usr/bin/
usr/bin/ntpctl
usr/bin/ntpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/ntp-units.d/
usr/lib/systemd/ntp-units.d/openntpd.list
usr/lib/systemd/system/
usr/lib/systemd/system/openntpd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/openntpd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/openntpd/
usr/share/licenses/openntpd/COPYING
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/ntpd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/ntpctl.8.gz
usr/share/man/man8/ntpd.8.gz
var/
var/db/
var/lib/
var/lib/ntp/
