# extracted/nrpe-4.1.1-1/files
%FILES%
etc/
etc/nrpe/
etc/nrpe/nrpe.cfg
usr/
usr/bin/
usr/bin/nrpe
usr/lib/
usr/lib/monitoring-plugins/
usr/lib/monitoring-plugins/check_nrpe
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nrpe.service
usr/lib/systemd/system/nrpe.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/nrpe.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nrpe.conf
usr/share/
usr/share/doc/
usr/share/doc/nrpe/
usr/share/doc/nrpe/README.SSL.md
usr/share/doc/nrpe/README.md
usr/share/doc/nrpe/SECURITY.md
usr/share/licenses/
usr/share/licenses/nrpe/
usr/share/licenses/nrpe/LEGAL
