# extracted/lxcfs-6.0.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/lxcfs
usr/lib/
usr/lib/lxcfs/
usr/lib/lxcfs/liblxcfs.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lxcfs.service
usr/share/
usr/share/lxc/
usr/share/lxc/config/
usr/share/lxc/config/common.conf.d/
usr/share/lxc/config/common.conf.d/00-lxcfs.conf
usr/share/lxcfs/
usr/share/lxcfs/lxc.mount.hook
usr/share/lxcfs/lxc.reboot.hook
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/lxcfs.1.gz
var/
var/lib/
var/lib/lxcfs/
