# extracted/wiresmith-0.4.4-1/files
%FILES%
usr/
usr/bin/
usr/bin/wiresmith
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wiresmith.service
usr/share/
usr/share/licenses/
usr/share/licenses/wiresmith/
usr/share/licenses/wiresmith/LICENSE
