# extracted/choria-io-0.23.0-1/files
%FILES%
etc/
etc/choria/
etc/choria/broker.conf
etc/choria/server.conf
etc/default/
etc/default/choria-server
etc/logrotate.d/
etc/logrotate.d/choria
usr/
usr/bin/
usr/bin/choria
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/choria-broker.service
usr/lib/systemd/system/choria-server.service
usr/share/
usr/share/doc/
usr/share/doc/choria-io/
usr/share/doc/choria-io/CHANGELOG.md
usr/share/doc/choria-io/CONFIGURATION.md
usr/share/doc/choria-io/Dockerfile
usr/share/doc/choria-io/NOTICE
usr/share/doc/choria-io/README.md
usr/share/licenses/
usr/share/licenses/choria-io/
usr/share/licenses/choria-io/LICENSE
