# extracted/mkinitcpio-systemd-tool-38-2/files
%FILES%
etc/
etc/mkinitcpio-systemd-tool/
etc/mkinitcpio-systemd-tool/config/
etc/mkinitcpio-systemd-tool/config/crypttab
etc/mkinitcpio-systemd-tool/config/fstab
etc/mkinitcpio-systemd-tool/config/initrd-nftables.conf
etc/mkinitcpio-systemd-tool/config/initrd-util-usb-hcd.conf
etc/mkinitcpio-systemd-tool/config/ntp.conf
etc/mkinitcpio-systemd-tool/mkinitcpio-systemd-tool.conf
etc/mkinitcpio-systemd-tool/network/
etc/mkinitcpio-systemd-tool/network/initrd-network.network
usr/
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/systemd-tool
usr/lib/initcpio/install/
usr/lib/initcpio/install/systemd-tool
usr/lib/mkinitcpio-systemd-tool/
usr/lib/mkinitcpio-systemd-tool/initrd-build.sh
usr/lib/mkinitcpio-systemd-tool/initrd-shell.sh
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/initrd-cryptsetup.path
usr/lib/systemd/system/initrd-cryptsetup.service
usr/lib/systemd/system/initrd-debug-progs.service
usr/lib/systemd/system/initrd-debug-shell.service
usr/lib/systemd/system/initrd-dropbear.service
usr/lib/systemd/system/initrd-emergency.service
usr/lib/systemd/system/initrd-emergency.target
usr/lib/systemd/system/initrd-network.service
usr/lib/systemd/system/initrd-nftables.service
usr/lib/systemd/system/initrd-ntpd.service
usr/lib/systemd/system/initrd-plymouth.path
usr/lib/systemd/system/initrd-plymouth.service
usr/lib/systemd/system/initrd-shell.service
usr/lib/systemd/system/initrd-sysroot-mount.service
usr/lib/systemd/system/initrd-tinysshd.service
usr/lib/systemd/system/initrd-util-pc-beep.service
usr/lib/systemd/system/initrd-util-usb-hcd.service
usr/share/
usr/share/doc/
usr/share/doc/mkinitcpio-systemd-tool/
usr/share/doc/mkinitcpio-systemd-tool/README.md
usr/share/licenses/
usr/share/licenses/mkinitcpio-systemd-tool/
usr/share/licenses/mkinitcpio-systemd-tool/LICENSE.md
