# extracted/opensearch-2.18.0-1/files
%FILES%
etc/
etc/default/
etc/default/opensearch
etc/opensearch/
etc/opensearch/jvm.options
etc/opensearch/jvm.options.d/
etc/opensearch/log4j2.properties
etc/opensearch/opensearch.yml
usr/
usr/bin/
usr/bin/opensearch
usr/bin/opensearch-keystore
usr/bin/opensearch-node
usr/bin/opensearch-plugin
usr/bin/opensearch-shard
usr/bin/opensearch-upgrade
usr/lib/
usr/lib/sysctl.d/
usr/lib/sysctl.d/opensearch.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opensearch-keystore.service
usr/lib/systemd/system/opensearch-keystore@.service
usr/lib/systemd/system/opensearch.service
usr/lib/systemd/system/opensearch@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/opensearch.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/opensearch.conf
usr/share/
usr/share/licenses/
usr/share/licenses/opensearch/
usr/share/licenses/opensearch/LICENSE.txt
usr/share/opensearch/
usr/share/opensearch/LICENSE.txt
usr/share/opensearch/NOTICE.txt
usr/share/opensearch/README.md
usr/share/opensearch/bin/
usr/share/opensearch/bin/opensearch
usr/share/opensearch/bin/opensearch-cli
usr/share/opensearch/bin/opensearch-env
usr/share/opensearch/bin/opensearch-env-from-file
usr/share/opensearch/bin/opensearch-keystore
usr/share/opensearch/bin/opensearch-node
usr/share/opensearch/bin/opensearch-plugin
usr/share/opensearch/bin/opensearch-shard
usr/share/opensearch/bin/opensearch-upgrade
usr/share/opensearch/config
usr/share/opensearch/data
usr/share/opensearch/lib/
usr/share/opensearch/lib/HdrHistogram-2.2.2.jar
usr/share/opensearch/lib/RoaringBitmap-1.3.0.jar
usr/share/opensearch/lib/jackson-core-2.17.2.jar
usr/share/opensearch/lib/jackson-dataformat-cbor-2.17.2.jar
usr/share/opensearch/lib/jackson-dataformat-smile-2.17.2.jar
usr/share/opensearch/lib/jackson-dataformat-yaml-2.17.2.jar
usr/share/opensearch/lib/jakarta.annotation-api-1.3.5.jar
usr/share/opensearch/lib/java-version-checker-2.18.0.jar
usr/share/opensearch/lib/jna-5.13.0.jar
usr/share/opensearch/lib/joda-time-2.12.7.jar
usr/share/opensearch/lib/jopt-simple-5.0.4.jar
usr/share/opensearch/lib/jts-core-1.15.0.jar
usr/share/opensearch/lib/jzlib-1.1.3.jar
usr/share/opensearch/lib/log4j-api-2.21.0.jar
usr/share/opensearch/lib/log4j-core-2.21.0.jar
usr/share/opensearch/lib/log4j-jul-2.21.0.jar
usr/share/opensearch/lib/lucene-analysis-common-9.12.0.jar
usr/share/opensearch/lib/lucene-backward-codecs-9.12.0.jar
usr/share/opensearch/lib/lucene-core-9.12.0.jar
usr/share/opensearch/lib/lucene-grouping-9.12.0.jar
usr/share/opensearch/lib/lucene-highlighter-9.12.0.jar
usr/share/opensearch/lib/lucene-join-9.12.0.jar
usr/share/opensearch/lib/lucene-memory-9.12.0.jar
usr/share/opensearch/lib/lucene-misc-9.12.0.jar
usr/share/opensearch/lib/lucene-queries-9.12.0.jar
usr/share/opensearch/lib/lucene-queryparser-9.12.0.jar
usr/share/opensearch/lib/lucene-sandbox-9.12.0.jar
usr/share/opensearch/lib/lucene-spatial-extras-9.12.0.jar
usr/share/opensearch/lib/lucene-spatial3d-9.12.0.jar
usr/share/opensearch/lib/lucene-suggest-9.12.0.jar
usr/share/opensearch/lib/opensearch-2.18.0.jar
usr/share/opensearch/lib/opensearch-cli-2.18.0.jar
usr/share/opensearch/lib/opensearch-common-2.18.0.jar
usr/share/opensearch/lib/opensearch-compress-2.18.0.jar
usr/share/opensearch/lib/opensearch-core-2.18.0.jar
usr/share/opensearch/lib/opensearch-geo-2.18.0.jar
usr/share/opensearch/lib/opensearch-launchers-2.18.0.jar
usr/share/opensearch/lib/opensearch-plugin-classloader-2.18.0.jar
usr/share/opensearch/lib/opensearch-secure-sm-2.18.0.jar
usr/share/opensearch/lib/opensearch-task-commons-2.18.0.jar
usr/share/opensearch/lib/opensearch-telemetry-2.18.0.jar
usr/share/opensearch/lib/opensearch-x-content-2.18.0.jar
usr/share/opensearch/lib/protobuf-java-3.25.5.jar
usr/share/opensearch/lib/reactive-streams-1.0.4.jar
usr/share/opensearch/lib/reactor-core-3.5.20.jar
usr/share/opensearch/lib/snakeyaml-2.1.jar
usr/share/opensearch/lib/spatial4j-0.7.jar
usr/share/opensearch/lib/t-digest-3.2.jar
usr/share/opensearch/lib/tools/
usr/share/opensearch/lib/tools/keystore-cli/
usr/share/opensearch/lib/tools/keystore-cli/keystore-cli-2.18.0.jar
usr/share/opensearch/lib/tools/plugin-cli/
usr/share/opensearch/lib/tools/plugin-cli/bc-fips-2.0.0.jar
usr/share/opensearch/lib/tools/plugin-cli/bcpg-fips-2.0.9.jar
usr/share/opensearch/lib/tools/plugin-cli/commons-compress-1.26.1.jar
usr/share/opensearch/lib/tools/plugin-cli/commons-io-2.16.0.jar
usr/share/opensearch/lib/tools/plugin-cli/opensearch-plugin-cli-2.18.0.jar
usr/share/opensearch/lib/tools/upgrade-cli/
usr/share/opensearch/lib/tools/upgrade-cli/jackson-annotations-2.17.2.jar
usr/share/opensearch/lib/tools/upgrade-cli/jackson-core-2.17.2.jar
usr/share/opensearch/lib/tools/upgrade-cli/jackson-databind-2.17.2.jar
usr/share/opensearch/lib/tools/upgrade-cli/opensearch-upgrade-cli-2.18.0.jar
usr/share/opensearch/lib/zstd-jni-1.5.5-5.jar
usr/share/opensearch/logs
usr/share/opensearch/modules/
usr/share/opensearch/modules/aggs-matrix-stats/
usr/share/opensearch/modules/aggs-matrix-stats/aggs-matrix-stats-client-2.18.0.jar
usr/share/opensearch/modules/aggs-matrix-stats/plugin-descriptor.properties
usr/share/opensearch/modules/analysis-common/
usr/share/opensearch/modules/analysis-common/analysis-common-2.18.0.jar
usr/share/opensearch/modules/analysis-common/plugin-descriptor.properties
usr/share/opensearch/modules/cache-common/
usr/share/opensearch/modules/cache-common/cache-common-2.18.0.jar
usr/share/opensearch/modules/cache-common/plugin-descriptor.properties
usr/share/opensearch/modules/cache-common/plugin-security.policy
usr/share/opensearch/modules/geo/
usr/share/opensearch/modules/geo/geo-2.18.0.jar
usr/share/opensearch/modules/geo/plugin-descriptor.properties
usr/share/opensearch/modules/ingest-common/
usr/share/opensearch/modules/ingest-common/ingest-common-2.18.0.jar
usr/share/opensearch/modules/ingest-common/jcodings-1.0.58.jar
usr/share/opensearch/modules/ingest-common/joni-2.2.1.jar
usr/share/opensearch/modules/ingest-common/opensearch-dissect-2.18.0.jar
usr/share/opensearch/modules/ingest-common/opensearch-grok-2.18.0.jar
usr/share/opensearch/modules/ingest-common/plugin-descriptor.properties
usr/share/opensearch/modules/ingest-geoip/
usr/share/opensearch/modules/ingest-geoip/GeoLite2-ASN.mmdb
usr/share/opensearch/modules/ingest-geoip/GeoLite2-City.mmdb
usr/share/opensearch/modules/ingest-geoip/GeoLite2-Country.mmdb
usr/share/opensearch/modules/ingest-geoip/geoip2-4.2.1.jar
usr/share/opensearch/modules/ingest-geoip/ingest-geoip-2.18.0.jar
usr/share/opensearch/modules/ingest-geoip/jackson-annotations-2.17.2.jar
usr/share/opensearch/modules/ingest-geoip/jackson-databind-2.17.2.jar
usr/share/opensearch/modules/ingest-geoip/maxmind-db-3.1.1.jar
usr/share/opensearch/modules/ingest-geoip/plugin-descriptor.properties
usr/share/opensearch/modules/ingest-geoip/plugin-security.policy
usr/share/opensearch/modules/ingest-user-agent/
usr/share/opensearch/modules/ingest-user-agent/ingest-user-agent-2.18.0.jar
usr/share/opensearch/modules/ingest-user-agent/plugin-descriptor.properties
usr/share/opensearch/modules/lang-expression/
usr/share/opensearch/modules/lang-expression/antlr4-runtime-4.13.1.jar
usr/share/opensearch/modules/lang-expression/asm-9.7.jar
usr/share/opensearch/modules/lang-expression/asm-commons-9.7.jar
usr/share/opensearch/modules/lang-expression/asm-tree-9.7.jar
usr/share/opensearch/modules/lang-expression/lang-expression-2.18.0.jar
usr/share/opensearch/modules/lang-expression/lucene-expressions-9.12.0.jar
usr/share/opensearch/modules/lang-expression/plugin-descriptor.properties
usr/share/opensearch/modules/lang-expression/plugin-security.policy
usr/share/opensearch/modules/lang-mustache/
usr/share/opensearch/modules/lang-mustache/compiler-0.9.14.jar
usr/share/opensearch/modules/lang-mustache/lang-mustache-client-2.18.0.jar
usr/share/opensearch/modules/lang-mustache/plugin-descriptor.properties
usr/share/opensearch/modules/lang-mustache/plugin-security.policy
usr/share/opensearch/modules/lang-painless/
usr/share/opensearch/modules/lang-painless/antlr4-runtime-4.13.1.jar
usr/share/opensearch/modules/lang-painless/asm-9.7.jar
usr/share/opensearch/modules/lang-painless/asm-analysis-9.7.jar
usr/share/opensearch/modules/lang-painless/asm-commons-9.7.jar
usr/share/opensearch/modules/lang-painless/asm-tree-9.7.jar
usr/share/opensearch/modules/lang-painless/asm-util-9.7.jar
usr/share/opensearch/modules/lang-painless/lang-painless-2.18.0.jar
usr/share/opensearch/modules/lang-painless/opensearch-scripting-painless-spi-2.18.0.jar
usr/share/opensearch/modules/lang-painless/plugin-descriptor.properties
usr/share/opensearch/modules/lang-painless/plugin-security.policy
usr/share/opensearch/modules/mapper-extras/
usr/share/opensearch/modules/mapper-extras/mapper-extras-client-2.18.0.jar
usr/share/opensearch/modules/mapper-extras/plugin-descriptor.properties
usr/share/opensearch/modules/opensearch-dashboards/
usr/share/opensearch/modules/opensearch-dashboards/commons-codec-1.16.1.jar
usr/share/opensearch/modules/opensearch-dashboards/commons-logging-1.2.jar
usr/share/opensearch/modules/opensearch-dashboards/httpasyncclient-4.1.5.jar
usr/share/opensearch/modules/opensearch-dashboards/httpclient-4.5.14.jar
usr/share/opensearch/modules/opensearch-dashboards/httpcore-4.4.16.jar
usr/share/opensearch/modules/opensearch-dashboards/httpcore-nio-4.4.16.jar
usr/share/opensearch/modules/opensearch-dashboards/opensearch-dashboards-2.18.0.jar
usr/share/opensearch/modules/opensearch-dashboards/opensearch-rest-client-2.18.0.jar
usr/share/opensearch/modules/opensearch-dashboards/opensearch-ssl-config-2.18.0.jar
usr/share/opensearch/modules/opensearch-dashboards/plugin-descriptor.properties
usr/share/opensearch/modules/opensearch-dashboards/reindex-client-2.18.0.jar
usr/share/opensearch/modules/parent-join/
usr/share/opensearch/modules/parent-join/parent-join-client-2.18.0.jar
usr/share/opensearch/modules/parent-join/plugin-descriptor.properties
usr/share/opensearch/modules/percolator/
usr/share/opensearch/modules/percolator/percolator-client-2.18.0.jar
usr/share/opensearch/modules/percolator/plugin-descriptor.properties
usr/share/opensearch/modules/rank-eval/
usr/share/opensearch/modules/rank-eval/plugin-descriptor.properties
usr/share/opensearch/modules/rank-eval/rank-eval-client-2.18.0.jar
usr/share/opensearch/modules/reindex/
usr/share/opensearch/modules/reindex/commons-codec-1.16.1.jar
usr/share/opensearch/modules/reindex/commons-logging-1.2.jar
usr/share/opensearch/modules/reindex/httpasyncclient-4.1.5.jar
usr/share/opensearch/modules/reindex/httpclient-4.5.14.jar
usr/share/opensearch/modules/reindex/httpcore-4.4.16.jar
usr/share/opensearch/modules/reindex/httpcore-nio-4.4.16.jar
usr/share/opensearch/modules/reindex/opensearch-rest-client-2.18.0.jar
usr/share/opensearch/modules/reindex/opensearch-ssl-config-2.18.0.jar
usr/share/opensearch/modules/reindex/parent-join/
usr/share/opensearch/modules/reindex/parent-join/plugin-descriptor.properties
usr/share/opensearch/modules/reindex/plugin-descriptor.properties
usr/share/opensearch/modules/reindex/plugin-security.policy
usr/share/opensearch/modules/reindex/reindex-client-2.18.0.jar
usr/share/opensearch/modules/reindex/transport-netty4/
usr/share/opensearch/modules/reindex/transport-netty4/plugin-descriptor.properties
usr/share/opensearch/modules/reindex/transport-netty4/plugin-security.policy
usr/share/opensearch/modules/repository-url/
usr/share/opensearch/modules/repository-url/plugin-descriptor.properties
usr/share/opensearch/modules/repository-url/plugin-security.policy
usr/share/opensearch/modules/repository-url/repository-url-2.18.0.jar
usr/share/opensearch/modules/search-pipeline-common/
usr/share/opensearch/modules/search-pipeline-common/plugin-descriptor.properties
usr/share/opensearch/modules/search-pipeline-common/search-pipeline-common-2.18.0.jar
usr/share/opensearch/modules/systemd/
usr/share/opensearch/modules/systemd/plugin-descriptor.properties
usr/share/opensearch/modules/systemd/plugin-security.policy
usr/share/opensearch/modules/systemd/systemd-2.18.0.jar
usr/share/opensearch/modules/transport-netty4/
usr/share/opensearch/modules/transport-netty4/netty-buffer-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-codec-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-codec-http-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-common-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-handler-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-resolver-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-transport-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/netty-transport-native-unix-common-4.1.114.Final.jar
usr/share/opensearch/modules/transport-netty4/plugin-descriptor.properties
usr/share/opensearch/modules/transport-netty4/plugin-security.policy
usr/share/opensearch/modules/transport-netty4/transport-netty4-client-2.18.0.jar
usr/share/opensearch/plugins/
var/
var/lib/
var/lib/opensearch/
var/log/
var/log/opensearch/
