# extracted/rp-pppoe-4.0-2/files
%FILES%
etc/
etc/ppp/
etc/ppp/pppoe-server-options
usr/
usr/bin/
usr/bin/pppoe
usr/bin/pppoe-relay
usr/bin/pppoe-server
usr/bin/pppoe-sniff
usr/lib/
usr/lib/rp-pppoe/
usr/lib/rp-pppoe/README
usr/lib/rp-pppoe/rp-pppoe.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/adsl.service
usr/share/
usr/share/doc/
usr/share/doc/rp-pppoe-4.0/
usr/share/doc/rp-pppoe-4.0/CHANGES
usr/share/doc/rp-pppoe-4.0/HOW-TO-CONNECT
usr/share/doc/rp-pppoe-4.0/KERNEL-MODE-PPPOE
usr/share/doc/rp-pppoe-4.0/LICENSE
usr/share/doc/rp-pppoe-4.0/README.md
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/pppoe-relay.8.gz
usr/share/man/man8/pppoe-server.8.gz
usr/share/man/man8/pppoe-sniff.8.gz
usr/share/man/man8/pppoe.8.gz
