# extracted/globalprotect-openconnect-1.4.9-2/files
%FILES%
etc/
etc/gpservice/
etc/gpservice/gp.conf
usr/
usr/bin/
usr/bin/gpclient
usr/bin/gpservice
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gpservice.service
usr/share/
usr/share/applications/
usr/share/applications/com.yuezk.qt.gpclient.desktop
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/com.yuezk.qt.GPService.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.yuezk.qt.GPService.conf
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/com.yuezk.qt.gpclient.svg
usr/share/metainfo/
usr/share/metainfo/com.yuezk.qt.gpclient.metainfo.xml
