# extracted/bumblebee-3.2.1-21/files
%FILES%
etc/
etc/bumblebee/
etc/bumblebee/bumblebee.conf
etc/bumblebee/xorg.conf.d/
etc/bumblebee/xorg.conf.d/10-dummy.conf
etc/bumblebee/xorg.conf.nouveau
etc/bumblebee/xorg.conf.nvidia
usr/
usr/bin/
usr/bin/bumblebee-bugreport
usr/bin/bumblebeed
usr/bin/optirun
usr/lib/
usr/lib/modprobe.d/
usr/lib/modprobe.d/bumblebee.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bumblebeed.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/bumblebee.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/99-bumblebee-nvidia-dev.rules
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/optirun
usr/share/doc/
usr/share/doc/bumblebee/
usr/share/doc/bumblebee/README.markdown
usr/share/doc/bumblebee/RELEASE_NOTES_3_2_1
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/bumblebeed.1.gz
usr/share/man/man1/optirun.1.gz
