# extracted/i2pd-2.54.0-1/files
%FILES%
etc/
etc/i2pd/
etc/i2pd/i2pd.conf
etc/i2pd/tunnels.conf
etc/i2pd/tunnels.d/
etc/logrotate.d/
etc/logrotate.d/i2pd
usr/
usr/bin/
usr/bin/i2pd
usr/include/
usr/include/i2pd/
usr/include/i2pd/AddressBook.h
usr/include/i2pd/BOB.h
usr/include/i2pd/Base.h
usr/include/i2pd/Blinding.h
usr/include/i2pd/CPU.h
usr/include/i2pd/ClientContext.h
usr/include/i2pd/Config.h
usr/include/i2pd/Crypto.h
usr/include/i2pd/CryptoKey.h
usr/include/i2pd/Datagram.h
usr/include/i2pd/Destination.h
usr/include/i2pd/ECIESX25519AEADRatchetSession.h
usr/include/i2pd/Ed25519.h
usr/include/i2pd/Elligator.h
usr/include/i2pd/FS.h
usr/include/i2pd/Family.h
usr/include/i2pd/Garlic.h
usr/include/i2pd/Gost.h
usr/include/i2pd/Gzip.h
usr/include/i2pd/HTTP.h
usr/include/i2pd/HTTPProxy.h
usr/include/i2pd/I18N.h
usr/include/i2pd/I18N_langs.h
usr/include/i2pd/I2CP.h
usr/include/i2pd/I2NPProtocol.h
usr/include/i2pd/I2PEndian.h
usr/include/i2pd/I2PService.h
usr/include/i2pd/I2PTunnel.h
usr/include/i2pd/Identity.h
usr/include/i2pd/KadDHT.h
usr/include/i2pd/LeaseSet.h
usr/include/i2pd/LittleBigEndian.h
usr/include/i2pd/Log.h
usr/include/i2pd/MatchedDestination.h
usr/include/i2pd/NTCP2.h
usr/include/i2pd/NetDbRequests.h
usr/include/i2pd/Profiling.h
usr/include/i2pd/Queue.h
usr/include/i2pd/Reseed.h
usr/include/i2pd/RouterContext.h
usr/include/i2pd/RouterInfo.h
usr/include/i2pd/SAM.h
usr/include/i2pd/SOCKS.h
usr/include/i2pd/SSU2.h
usr/include/i2pd/SSU2Session.h
usr/include/i2pd/Signature.h
usr/include/i2pd/Siphash.h
usr/include/i2pd/Socks5.h
usr/include/i2pd/Streaming.h
usr/include/i2pd/Tag.h
usr/include/i2pd/Timestamp.h
usr/include/i2pd/TransitTunnel.h
usr/include/i2pd/TransportSession.h
usr/include/i2pd/Transports.h
usr/include/i2pd/Tunnel.h
usr/include/i2pd/TunnelBase.h
usr/include/i2pd/TunnelConfig.h
usr/include/i2pd/TunnelEndpoint.h
usr/include/i2pd/TunnelGateway.h
usr/include/i2pd/TunnelPool.h
usr/include/i2pd/UDPTunnel.h
usr/include/i2pd/api.h
usr/include/i2pd/util.h
usr/include/i2pd/version.h
usr/lib/
usr/lib/libi2pd.so
usr/lib/libi2pdclient.so
usr/lib/libi2pdlang.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/i2pd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/i2pd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/i2pd.conf
usr/share/
usr/share/doc/
usr/share/doc/i2pd/
usr/share/doc/i2pd/tunnels.d/
usr/share/doc/i2pd/tunnels.d/IRC-Ilita.conf
usr/share/doc/i2pd/tunnels.d/IRC-Irc2P.conf
usr/share/doc/i2pd/tunnels.d/README
usr/share/i2pd/
usr/share/i2pd/certificates/
usr/share/i2pd/certificates/family/
usr/share/i2pd/certificates/family/gostcoin.crt
usr/share/i2pd/certificates/family/i2p-dev.crt
usr/share/i2pd/certificates/family/i2pd-dev.crt
usr/share/i2pd/certificates/family/mca2-i2p.crt
usr/share/i2pd/certificates/family/stormycloud.crt
usr/share/i2pd/certificates/family/volatile.crt
usr/share/i2pd/certificates/reseed/
usr/share/i2pd/certificates/reseed/acetone_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/admin_at_stormycloud.org.crt
usr/share/i2pd/certificates/reseed/arnavbhatt288_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/creativecowpat_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/echelon3_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/hankhill19580_at_gmail.com.crt
usr/share/i2pd/certificates/reseed/hiduser0_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/hottuna_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/i2p-reseed_at_mk16.de.crt
usr/share/i2pd/certificates/reseed/igor_at_novg.net.crt
usr/share/i2pd/certificates/reseed/lazygravy_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/orignal_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/r4sas-reseed_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/rambler_at_mail.i2p.crt
usr/share/i2pd/certificates/reseed/reseed_at_diva.exchange.crt
usr/share/licenses/
usr/share/licenses/i2pd/
usr/share/licenses/i2pd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/i2pd.1.gz
