# extracted/woodpecker-agent-2.7.1-1/files
%FILES%
etc/
etc/woodpecker/
etc/woodpecker/agent.env
usr/
usr/bin/
usr/bin/woodpecker-agent
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/woodpecker-agent.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/woodpecker-agent.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/woodpecker-agent.conf
