# extracted/umurmur-0.2.20-3/files
%FILES%
etc/
etc/umurmur/
etc/umurmur/umurmur.conf
usr/
usr/bin/
usr/bin/umurmurd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/umurmur.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/umurmur.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/umurmur.conf
usr/share/
usr/share/doc/
usr/share/doc/umurmur/
usr/share/doc/umurmur/AUTHORS
usr/share/doc/umurmur/ChangeLog
usr/share/doc/umurmur/README.md
usr/share/licenses/
usr/share/licenses/umurmur/
usr/share/licenses/umurmur/LICENSE
