# extracted/pacoloco-1.7-1/files
%FILES%
etc/
etc/pacoloco.yaml
usr/
usr/bin/
usr/bin/pacoloco
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pacoloco.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pacoloco.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pacoloco.conf
