# extracted/acme-redirect-0.7.0-1/files
%FILES%
etc/
etc/acme-redirect.conf
etc/acme-redirect.d/
etc/acme-redirect.d/example.com.conf.sample
usr/
usr/bin/
usr/bin/acme-redirect
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/acme-redirect-renew.service
usr/lib/systemd/system/acme-redirect-renew.timer
usr/lib/systemd/system/acme-redirect.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/acme-redirect.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/acme-redirect.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/acme-redirect
usr/share/doc/
usr/share/doc/acme-redirect/
usr/share/doc/acme-redirect/README.md
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/acme-redirect.fish
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/acme-redirect.1.gz
usr/share/man/man5/
usr/share/man/man5/acme-redirect.conf.5.gz
usr/share/man/man5/acme-redirect.d.5.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_acme-redirect
