# extracted/libldm-0.2.5-3/files
%FILES%
usr/
usr/bin/
usr/bin/ldmtool
usr/include/
usr/include/ldm.h
usr/lib/
usr/lib/libldm-1.0.so
usr/lib/libldm-1.0.so.0
usr/lib/libldm-1.0.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/ldm-1.0.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ldmtool.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ldmtool.1.gz
