# extracted/cozy-stack-1:1.6.33-1/files
%FILES%
etc/
etc/cozy/
etc/cozy/cozy.yml
usr/
usr/bin/
usr/bin/cozy-stack
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cozy-stack.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/cozy-stack.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/cozy-stack.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/cozy-stack
usr/share/cozy/
usr/share/cozy/cozy.example.yaml
usr/share/cozy/konnector-node-run.sh
usr/share/cozy/konnector-nsjail-run.sh
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/cozy-stack.fish
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_cozy-stack
