# extracted/virtualbox-guest-utils-7.1.4-2/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/vboxclient.desktop
usr/
usr/bin/
usr/bin/VBoxAudioTest
usr/bin/VBoxClient
usr/bin/VBoxClient-all
usr/bin/VBoxControl
usr/bin/VBoxDRMClient
usr/bin/VBoxService
usr/bin/vboxwl
usr/lib/
usr/lib/security/
usr/lib/security/pam_vbox.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vboxdrmclient.path
usr/lib/systemd/system/vboxdrmclient.service
usr/lib/systemd/system/vboxservice.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/virtualbox-guest-utils.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-vboxguest.rules
usr/share/
usr/share/licenses/
usr/share/licenses/virtualbox-guest-utils/
usr/share/licenses/virtualbox-guest-utils/LICENSE
usr/share/licenses/virtualbox-guest-utils/LICENSE.CDDL
