# extracted/fcron-3.3.1-2/files
%FILES%
etc/
etc/cron.daily/
etc/cron.hourly/
etc/cron.monthly/
etc/cron.weekly/
etc/fcron/
etc/fcron/fcron.allow
etc/fcron/fcron.conf
etc/fcron/fcron.deny
etc/fcron/pam.conf
etc/pam.d/
etc/pam.d/fcron
etc/pam.d/fcrontab
usr/
usr/bin/
usr/bin/fcron
usr/bin/fcrondyn
usr/bin/fcronsighup
usr/bin/fcrontab
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/fcron.service
usr/lib/systemd/system/fcronreboot.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/fcron.conf
usr/share/
usr/share/man/
usr/share/man/fr/
usr/share/man/fr/man1/
usr/share/man/fr/man1/fcrondyn.1.gz
usr/share/man/fr/man1/fcrontab.1.gz
usr/share/man/fr/man5/
usr/share/man/fr/man5/fcron.conf.5.gz
usr/share/man/fr/man5/fcrontab.5.gz
usr/share/man/fr/man8/
usr/share/man/fr/man8/fcron.8.gz
usr/share/man/man1/
usr/share/man/man1/fcrondyn.1.gz
usr/share/man/man1/fcrontab.1.gz
usr/share/man/man5/
usr/share/man/man5/fcron.conf.5.gz
usr/share/man/man5/fcrontab.5.gz
usr/share/man/man8/
usr/share/man/man8/fcron.8.gz
var/
var/spool/
var/spool/fcron/
var/spool/fcron/systab.orig
