# extracted/cri-o-1.31.3-1/files
%FILES%
etc/
etc/cni/
etc/cni/net.d/
etc/crictl.yaml
etc/crio/
etc/crio/crio.conf
etc/crio/crio.conf.d/
usr/
usr/bin/
usr/bin/crio
usr/bin/pinns
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/cri-o.conf
usr/lib/sysctl.d/
usr/lib/sysctl.d/90-cri-o.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/crio-wipe.service
usr/lib/systemd/system/crio.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/crio
usr/share/containers/
usr/share/containers/oci/
usr/share/containers/oci/hooks.d/
usr/share/doc/
usr/share/doc/cri-o/
usr/share/doc/cri-o/CONTRIBUTING.md
usr/share/doc/cri-o/README.md
usr/share/doc/cri-o/awesome.md
usr/share/doc/cri-o/code-of-conduct.md
usr/share/doc/cri-o/examples/
usr/share/doc/cri-o/examples/cni/
usr/share/doc/cri-o/examples/cni/10-crio-bridge.conflist
usr/share/doc/cri-o/examples/cni/11-crio-ipv4-bridge.conflist
usr/share/doc/cri-o/examples/cni/99-loopback.conflist
usr/share/doc/cri-o/examples/cni/README.md
usr/share/doc/cri-o/examples/metrics-exporter/
usr/share/doc/cri-o/examples/metrics-exporter/Containerfile
usr/share/doc/cri-o/examples/metrics-exporter/cluster.yaml
usr/share/doc/cri-o/examples/metrics-exporter/dashboard.gif
usr/share/doc/cri-o/examples/metrics-exporter/dashboard.json
usr/share/doc/cri-o/examples/metrics-exporter/main.go
usr/share/doc/cri-o/install.md
usr/share/doc/cri-o/transfer.md
usr/share/doc/cri-o/tutorial.md
usr/share/doc/cri-o/tutorials/
usr/share/doc/cri-o/tutorials/crictl.md
usr/share/doc/cri-o/tutorials/crio-in-kind.md
usr/share/doc/cri-o/tutorials/debugging.md
usr/share/doc/cri-o/tutorials/decryption.md
usr/share/doc/cri-o/tutorials/kubeadm.md
usr/share/doc/cri-o/tutorials/kubernetes.md
usr/share/doc/cri-o/tutorials/metrics.md
usr/share/doc/cri-o/tutorials/tracing.md
usr/share/doc/cri-o/tutorials/userns.md
usr/share/fish/
usr/share/fish/completions/
usr/share/fish/completions/crio.fish
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/crio.conf.5.gz
usr/share/man/man5/crio.conf.d.5.gz
usr/share/man/man8/
usr/share/man/man8/crio.8.gz
usr/share/oci-umount/
usr/share/oci-umount/oci-umount.d/
usr/share/oci-umount/oci-umount.d/crio-umount.conf
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_crio
