# extracted/kube-proxy-1.31.3-1/files
%FILES%
etc/
etc/kubernetes/
etc/kubernetes/kube-proxy.env
usr/
usr/bin/
usr/bin/kube-proxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/kube-proxy.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/kube-proxy.1.gz
var/
var/lib/
var/lib/kube-proxy/
