# extracted/nsd-4.11.0-1/files
%FILES%
etc/
etc/nsd/
etc/nsd/nsd.conf.sample
usr/
usr/bin/
usr/bin/nsd
usr/bin/nsd-checkconf
usr/bin/nsd-checkzone
usr/bin/nsd-control
usr/bin/nsd-control-setup
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nsd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nsd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nsd.conf
usr/share/
usr/share/doc/
usr/share/doc/nsd/
usr/share/doc/nsd/CREDITS
usr/share/doc/nsd/ChangeLog
usr/share/doc/nsd/LICENSE
usr/share/doc/nsd/NEW-CFG-OPTION
usr/share/doc/nsd/NSD-4-features
usr/share/doc/nsd/NSD-DATABASE
usr/share/doc/nsd/NSD-DIFFFILE
usr/share/doc/nsd/NSD-FOR-BIND-USERS
usr/share/doc/nsd/NSD-VERIFY-MODS
usr/share/doc/nsd/README
usr/share/doc/nsd/README.icc
usr/share/doc/nsd/README.svn
usr/share/doc/nsd/RELNOTES
usr/share/doc/nsd/REQUIREMENTS
usr/share/doc/nsd/TESTPLAN
usr/share/doc/nsd/TODO
usr/share/doc/nsd/UPGRADING
usr/share/doc/nsd/coding-style
usr/share/doc/nsd/differences.tex
usr/share/doc/nsd/manual/
usr/share/doc/nsd/manual/catalog-zones.rst
usr/share/doc/nsd/manual/conf.py.in
usr/share/doc/nsd/manual/configuration.rst
usr/share/doc/nsd/manual/index.rst
usr/share/doc/nsd/manual/installation.rst
usr/share/doc/nsd/manual/manpages/
usr/share/doc/nsd/manual/manpages/nsd-checkconf.rst
usr/share/doc/nsd/manual/manpages/nsd-checkzone.rst
usr/share/doc/nsd/manual/manpages/nsd-control.rst
usr/share/doc/nsd/manual/manpages/nsd.conf.rst
usr/share/doc/nsd/manual/manpages/nsd.rst
usr/share/doc/nsd/manual/reference/
usr/share/doc/nsd/manual/reference/configure-options.rst
usr/share/doc/nsd/manual/reference/log-diagnosis.rst
usr/share/doc/nsd/manual/requirements.txt
usr/share/doc/nsd/manual/resources/
usr/share/doc/nsd/manual/resources/css/
usr/share/doc/nsd/manual/resources/css/dark.css
usr/share/doc/nsd/manual/resources/css/light.css
usr/share/doc/nsd/manual/resources/favicon.ico
usr/share/doc/nsd/manual/resources/nsd-duotone-white.png
usr/share/doc/nsd/manual/resources/nsd-duotone-white.svg
usr/share/doc/nsd/manual/running/
usr/share/doc/nsd/manual/running/interfaces.rst
usr/share/doc/nsd/manual/running/logging.rst
usr/share/doc/nsd/manual/running/tuning.rst
usr/share/doc/nsd/manual/running/using-tsig.rst
usr/share/doc/nsd/manual/running/zone-expiry.rst
usr/share/doc/nsd/manual/zonefile.rst
usr/share/licenses/
usr/share/licenses/nsd/
usr/share/licenses/nsd/LICENSE
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/nsd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/nsd-checkconf.8.gz
usr/share/man/man8/nsd-checkzone.8.gz
usr/share/man/man8/nsd-control.8.gz
usr/share/man/man8/nsd.8.gz
var/
var/db/
var/db/nsd/
