# extracted/systemd-bootchart-235-2/files
%FILES%
etc/
etc/systemd/
etc/systemd/bootchart.conf
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/systemd-bootchart.service
usr/lib/systemd/systemd-bootchart
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/systemd-bootchart.1.gz
usr/share/man/man5/
usr/share/man/man5/bootchart.conf.5.gz
usr/share/man/man5/bootchart.conf.d.5.gz
