# extracted/consul-template-0.39.1-1/files
%FILES%
etc/
etc/consul-template/
etc/consul-template/consul-template.hcl
usr/
usr/bin/
usr/bin/consul-template
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/consul-template.service
usr/share/
usr/share/doc/
usr/share/doc/consul-template/
usr/share/doc/consul-template/README.md
usr/share/doc/consul-template/configuration.md
usr/share/doc/consul-template/examples/
usr/share/doc/consul-template/examples/apache.md
usr/share/doc/consul-template/examples/haproxy-connect-proxy.md
usr/share/doc/consul-template/examples/haproxy-connect-proxy/
usr/share/doc/consul-template/examples/haproxy-connect-proxy/run-haproxy-connect-proxy
usr/share/doc/consul-template/examples/haproxy.md
usr/share/doc/consul-template/examples/join.md
usr/share/doc/consul-template/examples/nginx-connect-proxy.md
usr/share/doc/consul-template/examples/nginx-connect-proxy/
usr/share/doc/consul-template/examples/nginx-connect-proxy/run-nginx-connect-proxy
usr/share/doc/consul-template/examples/nginx.md
usr/share/doc/consul-template/examples/services.md
usr/share/doc/consul-template/examples/varnish.md
usr/share/doc/consul-template/examples/vault-pki.md
usr/share/doc/consul-template/examples/vault-transit.md
usr/share/doc/consul-template/modes.md
usr/share/doc/consul-template/observability.md
usr/share/doc/consul-template/plugins.md
usr/share/doc/consul-template/templating-language.md
usr/share/licenses/
usr/share/licenses/consul-template/
usr/share/licenses/consul-template/LICENSE
