# extracted/notmuch-mutt-0.38.3-4/files
%FILES%
etc/
etc/Muttrc.d/
etc/Muttrc.d/notmuch-mutt.rc
usr/
usr/bin/
usr/bin/notmuch-mutt
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/notmuch-mutt.1.gz
