# extracted/incus-6.7-1/files
%FILES%
usr/
usr/bin/
usr/bin/incus
usr/bin/incus-agent
usr/bin/incus-user
usr/bin/incusd
usr/bin/lxd-to-incus
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/incus-user.service
usr/lib/systemd/system/incus-user.socket
usr/lib/systemd/system/incus.service
usr/lib/systemd/system/incus.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/incus.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/incus
usr/share/doc/
usr/share/doc/incus/
usr/share/doc/incus/README.md
usr/share/doc/incus/api-extensions.md
usr/share/doc/incus/api.md
usr/share/doc/incus/architectures.md
usr/share/doc/incus/authentication.md
usr/share/doc/incus/authorization.md
usr/share/doc/incus/backup.md
usr/share/doc/incus/client.md
usr/share/doc/incus/cloud-init.md
usr/share/doc/incus/clustering.md
usr/share/doc/incus/conf.py
usr/share/doc/incus/config-options.md
usr/share/doc/incus/config_options.txt
usr/share/doc/incus/config_options_cheat_sheet.md
usr/share/doc/incus/container-environment.md
usr/share/doc/incus/contributing.md
usr/share/doc/incus/daemon-behavior.md
usr/share/doc/incus/database.md
usr/share/doc/incus/debugging.md
usr/share/doc/incus/dev-incus.md
usr/share/doc/incus/doc-cheat-sheet.md
usr/share/doc/incus/environment.md
usr/share/doc/incus/events.md
usr/share/doc/incus/explanation/
usr/share/doc/incus/explanation/clustering.md
usr/share/doc/incus/explanation/containers_and_vms.md
usr/share/doc/incus/explanation/instance_config.md
usr/share/doc/incus/explanation/instances.md
usr/share/doc/incus/explanation/networks.md
usr/share/doc/incus/explanation/performance_tuning.md
usr/share/doc/incus/explanation/projects.md
usr/share/doc/incus/explanation/security.md
usr/share/doc/incus/explanation/storage.md
usr/share/doc/incus/external_resources.md
usr/share/doc/incus/faq.md
usr/share/doc/incus/general.md
usr/share/doc/incus/howto/
usr/share/doc/incus/howto/benchmark_performance.md
usr/share/doc/incus/howto/cluster_access.md
usr/share/doc/incus/howto/cluster_config_networks.md
usr/share/doc/incus/howto/cluster_config_storage.md
usr/share/doc/incus/howto/cluster_form.md
usr/share/doc/incus/howto/cluster_groups.md
usr/share/doc/incus/howto/cluster_manage.md
usr/share/doc/incus/howto/cluster_manage_instance.md
usr/share/doc/incus/howto/cluster_recover.md
usr/share/doc/incus/howto/disaster_recovery.md
usr/share/doc/incus/howto/images_copy.md
usr/share/doc/incus/howto/images_create.md
usr/share/doc/incus/howto/images_manage.md
usr/share/doc/incus/howto/images_profiles.md
usr/share/doc/incus/howto/images_remote.md
usr/share/doc/incus/howto/import_machines_to_instances.md
usr/share/doc/incus/howto/incus_alias.md
usr/share/doc/incus/howto/initialize.md
usr/share/doc/incus/howto/instances_access_files.md
usr/share/doc/incus/howto/instances_backup.md
usr/share/doc/incus/howto/instances_configure.md
usr/share/doc/incus/howto/instances_console.md
usr/share/doc/incus/howto/instances_create.md
usr/share/doc/incus/howto/instances_manage.md
usr/share/doc/incus/howto/instances_routed_nic_vm.md
usr/share/doc/incus/howto/instances_troubleshoot.md
usr/share/doc/incus/howto/migrate_from_lxc.md
usr/share/doc/incus/howto/move_instances.md
usr/share/doc/incus/howto/network_acls.md
usr/share/doc/incus/howto/network_bgp.md
usr/share/doc/incus/howto/network_bridge_firewalld.md
usr/share/doc/incus/howto/network_bridge_resolved.md
usr/share/doc/incus/howto/network_configure.md
usr/share/doc/incus/howto/network_create.md
usr/share/doc/incus/howto/network_forwards.md
usr/share/doc/incus/howto/network_increase_bandwidth.md
usr/share/doc/incus/howto/network_integrations.md
usr/share/doc/incus/howto/network_ipam.md
usr/share/doc/incus/howto/network_load_balancers.md
usr/share/doc/incus/howto/network_ovn_peers.md
usr/share/doc/incus/howto/network_ovn_setup.md
usr/share/doc/incus/howto/network_zones.md
usr/share/doc/incus/howto/projects_confine.md
usr/share/doc/incus/howto/projects_create.md
usr/share/doc/incus/howto/projects_work.md
usr/share/doc/incus/howto/server_configure.md
usr/share/doc/incus/howto/server_expose.md
usr/share/doc/incus/howto/server_migrate_lxd.md
usr/share/doc/incus/howto/storage_backup_volume.md
usr/share/doc/incus/howto/storage_buckets.md
usr/share/doc/incus/howto/storage_create_instance.md
usr/share/doc/incus/howto/storage_move_volume.md
usr/share/doc/incus/howto/storage_pools.md
usr/share/doc/incus/howto/storage_volumes.md
usr/share/doc/incus/image-handling.md
usr/share/doc/incus/images.md
usr/share/doc/incus/images/
usr/share/doc/incus/images/UI/
usr/share/doc/incus/images/UI/limits_memory_example.png
usr/share/doc/incus/images/application-vs-system-containers.svg
usr/share/doc/incus/images/grafana_add_datasource.png
usr/share/doc/incus/images/grafana_configure_datasource.png
usr/share/doc/incus/images/grafana_dashboard_id.png
usr/share/doc/incus/images/grafana_dashboard_import.png
usr/share/doc/incus/images/grafana_dashboard_select_datasource.png
usr/share/doc/incus/images/grafana_instances.png
usr/share/doc/incus/images/grafana_resources.png
usr/share/doc/incus/images/grafana_select_prometheus.png
usr/share/doc/incus/images/ui_console.png
usr/share/doc/incus/images/ui_security_warning.png
usr/share/doc/incus/images/ui_set_up_certificates.png
usr/share/doc/incus/images/virtual-machines-vs-system-containers.svg
usr/share/doc/incus/index.md
usr/share/doc/incus/installing.md
usr/share/doc/incus/instance-exec.md
usr/share/doc/incus/instances.md
usr/share/doc/incus/internals.md
usr/share/doc/incus/metrics.md
usr/share/doc/incus/migration.md
usr/share/doc/incus/networks.md
usr/share/doc/incus/packaging.md
usr/share/doc/incus/profiles.md
usr/share/doc/incus/projects.md
usr/share/doc/incus/reference/
usr/share/doc/incus/reference/cluster_member_config.md
usr/share/doc/incus/reference/devices.md
usr/share/doc/incus/reference/devices_disk.md
usr/share/doc/incus/reference/devices_gpu.md
usr/share/doc/incus/reference/devices_infiniband.md
usr/share/doc/incus/reference/devices_nic.md
usr/share/doc/incus/reference/devices_none.md
usr/share/doc/incus/reference/devices_pci.md
usr/share/doc/incus/reference/devices_proxy.md
usr/share/doc/incus/reference/devices_tpm.md
usr/share/doc/incus/reference/devices_unix_block.md
usr/share/doc/incus/reference/devices_unix_char.md
usr/share/doc/incus/reference/devices_unix_hotplug.md
usr/share/doc/incus/reference/devices_usb.md
usr/share/doc/incus/reference/image_format.md
usr/share/doc/incus/reference/image_servers.md
usr/share/doc/incus/reference/instance_options.md
usr/share/doc/incus/reference/instance_properties.md
usr/share/doc/incus/reference/instance_units.md
usr/share/doc/incus/reference/manpages.md
usr/share/doc/incus/reference/manpages/
usr/share/doc/incus/reference/manpages/readme.txt
usr/share/doc/incus/reference/network_bridge.md
usr/share/doc/incus/reference/network_external.md
usr/share/doc/incus/reference/network_macvlan.md
usr/share/doc/incus/reference/network_ovn.md
usr/share/doc/incus/reference/network_physical.md
usr/share/doc/incus/reference/network_sriov.md
usr/share/doc/incus/reference/projects.md
usr/share/doc/incus/reference/provided_metrics.md
usr/share/doc/incus/reference/server_settings.md
usr/share/doc/incus/reference/standard_devices.md
usr/share/doc/incus/reference/storage_btrfs.md
usr/share/doc/incus/reference/storage_ceph.md
usr/share/doc/incus/reference/storage_cephfs.md
usr/share/doc/incus/reference/storage_cephobject.md
usr/share/doc/incus/reference/storage_dir.md
usr/share/doc/incus/reference/storage_drivers.md
usr/share/doc/incus/reference/storage_lvm.md
usr/share/doc/incus/reference/storage_zfs.md
usr/share/doc/incus/remotes.md
usr/share/doc/incus/requirements.md
usr/share/doc/incus/rest-api-spec.md
usr/share/doc/incus/rest-api.md
usr/share/doc/incus/rest-api.yaml
usr/share/doc/incus/security.md
usr/share/doc/incus/server.md
usr/share/doc/incus/server_config.md
usr/share/doc/incus/storage.md
usr/share/doc/incus/substitutions.yaml
usr/share/doc/incus/support.md
usr/share/doc/incus/syscall-interception.md
usr/share/doc/incus/third_party.md
usr/share/doc/incus/tutorial/
usr/share/doc/incus/tutorial/first_steps.md
usr/share/doc/incus/userns-idmap.md
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/incus.fish
usr/share/licenses/
usr/share/licenses/incus/
usr/share/licenses/incus/LICENCE
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_incus
var/
var/log/
var/log/incus/
