# extracted/alertmanager-0.27.0-2/files
%FILES%
etc/
etc/alertmanager/
etc/alertmanager/alertmanager.yml
etc/alertmanager/template/
etc/conf.d/
etc/conf.d/alertmanager
usr/
usr/bin/
usr/bin/alertmanager
usr/bin/amtool
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/alertmanager.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/alertmanager.conf
var/
var/lib/
var/lib/alertmanager/
