# extracted/throttled-0.10.0-4/files
%FILES%
etc/
etc/throttled.conf
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/throttled.service
usr/lib/throttled/
usr/lib/throttled/__pycache__/
usr/lib/throttled/__pycache__/mmio.cpython-313.pyc
usr/lib/throttled/__pycache__/throttled.cpython-313.pyc
usr/lib/throttled/mmio.py
usr/lib/throttled/throttled.py
usr/share/
usr/share/licenses/
usr/share/licenses/throttled/
usr/share/licenses/throttled/LICENSE
