# extracted/syncthing-discosrv-1.28.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/syncthing-discosrv
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/syncthing-discosrv.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/syncthing-discosrv.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/syncthing-discosrv.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/syncthing-discosrv.1.gz
