# extracted/tarantool-2.11.1-3/files
%FILES%
etc/
etc/default/
etc/default/tarantool
etc/logrotate.d/
etc/logrotate.d/tarantool
etc/tarantool/
etc/tarantool/instances.available/
etc/tarantool/instances.available/example.lua
usr/
usr/bin/
usr/bin/tarantool
usr/bin/tarantoolctl
usr/include/
usr/include/tarantool/
usr/include/tarantool/curl/
usr/include/tarantool/curl/curl.h
usr/include/tarantool/curl/curlver.h
usr/include/tarantool/curl/easy.h
usr/include/tarantool/curl/header.h
usr/include/tarantool/curl/mprintf.h
usr/include/tarantool/curl/multi.h
usr/include/tarantool/curl/options.h
usr/include/tarantool/curl/stdcheaders.h
usr/include/tarantool/curl/system.h
usr/include/tarantool/curl/typecheck-gcc.h
usr/include/tarantool/curl/urlapi.h
usr/include/tarantool/curl/websockets.h
usr/include/tarantool/lauxlib.h
usr/include/tarantool/lmisclib.h
usr/include/tarantool/lua.h
usr/include/tarantool/lua.hpp
usr/include/tarantool/luaconf.h
usr/include/tarantool/luajit.h
usr/include/tarantool/lualib.h
usr/include/tarantool/module.h
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tarantool@.service
usr/lib/systemd/system/tarantool@.service.d/
usr/lib/systemd/system/tarantool@.service.d/override.conf
usr/lib/sysusers.d/
usr/lib/sysusers.d/tarantool.conf
usr/lib/tarantool/
usr/lib/tmpfiles.d/
usr/share/
usr/share/licenses/
usr/share/licenses/tarantool/
usr/share/licenses/tarantool/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/tarantool.1.gz
usr/share/man/man1/tarantoolctl.1.gz
usr/share/tarantool/
usr/share/tarantool/luarocks/
usr/share/tarantool/luarocks/admin/
usr/share/tarantool/luarocks/admin/cache.lua
usr/share/tarantool/luarocks/admin/cmd/
usr/share/tarantool/luarocks/admin/cmd/add.lua
usr/share/tarantool/luarocks/admin/cmd/make_manifest.lua
usr/share/tarantool/luarocks/admin/cmd/refresh_cache.lua
usr/share/tarantool/luarocks/admin/cmd/remove.lua
usr/share/tarantool/luarocks/admin/index.lua
usr/share/tarantool/luarocks/build.lua
usr/share/tarantool/luarocks/build/
usr/share/tarantool/luarocks/build/builtin.lua
usr/share/tarantool/luarocks/build/cmake.lua
usr/share/tarantool/luarocks/build/command.lua
usr/share/tarantool/luarocks/build/make.lua
usr/share/tarantool/luarocks/cmd.lua
usr/share/tarantool/luarocks/cmd/
usr/share/tarantool/luarocks/cmd/build.lua
usr/share/tarantool/luarocks/cmd/config.lua
usr/share/tarantool/luarocks/cmd/doc.lua
usr/share/tarantool/luarocks/cmd/download.lua
usr/share/tarantool/luarocks/cmd/help.lua
usr/share/tarantool/luarocks/cmd/init.lua
usr/share/tarantool/luarocks/cmd/install.lua
usr/share/tarantool/luarocks/cmd/lint.lua
usr/share/tarantool/luarocks/cmd/list.lua
usr/share/tarantool/luarocks/cmd/make.lua
usr/share/tarantool/luarocks/cmd/new_version.lua
usr/share/tarantool/luarocks/cmd/pack.lua
usr/share/tarantool/luarocks/cmd/path.lua
usr/share/tarantool/luarocks/cmd/purge.lua
usr/share/tarantool/luarocks/cmd/remove.lua
usr/share/tarantool/luarocks/cmd/search.lua
usr/share/tarantool/luarocks/cmd/show.lua
usr/share/tarantool/luarocks/cmd/test.lua
usr/share/tarantool/luarocks/cmd/unpack.lua
usr/share/tarantool/luarocks/cmd/upload.lua
usr/share/tarantool/luarocks/cmd/which.lua
usr/share/tarantool/luarocks/cmd/write_rockspec.lua
usr/share/tarantool/luarocks/core/
usr/share/tarantool/luarocks/core/cfg.lua
usr/share/tarantool/luarocks/core/dir.lua
usr/share/tarantool/luarocks/core/hardcoded.lua
usr/share/tarantool/luarocks/core/manif.lua
usr/share/tarantool/luarocks/core/path.lua
usr/share/tarantool/luarocks/core/persist.lua
usr/share/tarantool/luarocks/core/sysdetect.lua
usr/share/tarantool/luarocks/core/util.lua
usr/share/tarantool/luarocks/core/vers.lua
usr/share/tarantool/luarocks/deps.lua
usr/share/tarantool/luarocks/dir.lua
usr/share/tarantool/luarocks/download.lua
usr/share/tarantool/luarocks/fetch.lua
usr/share/tarantool/luarocks/fetch/
usr/share/tarantool/luarocks/fetch/cvs.lua
usr/share/tarantool/luarocks/fetch/git.lua
usr/share/tarantool/luarocks/fetch/git_file.lua
usr/share/tarantool/luarocks/fetch/git_http.lua
usr/share/tarantool/luarocks/fetch/git_https.lua
usr/share/tarantool/luarocks/fetch/git_ssh.lua
usr/share/tarantool/luarocks/fetch/hg.lua
usr/share/tarantool/luarocks/fetch/hg_http.lua
usr/share/tarantool/luarocks/fetch/hg_https.lua
usr/share/tarantool/luarocks/fetch/hg_ssh.lua
usr/share/tarantool/luarocks/fetch/sscm.lua
usr/share/tarantool/luarocks/fetch/svn.lua
usr/share/tarantool/luarocks/fs.lua
usr/share/tarantool/luarocks/fs/
usr/share/tarantool/luarocks/fs/freebsd.lua
usr/share/tarantool/luarocks/fs/lua.lua
usr/share/tarantool/luarocks/fs/tools.lua
usr/share/tarantool/luarocks/fs/unix.lua
usr/share/tarantool/luarocks/fs/unix/
usr/share/tarantool/luarocks/fs/unix/tools.lua
usr/share/tarantool/luarocks/fs/win32.lua
usr/share/tarantool/luarocks/fs/win32/
usr/share/tarantool/luarocks/fs/win32/tools.lua
usr/share/tarantool/luarocks/fun.lua
usr/share/tarantool/luarocks/loader.lua
usr/share/tarantool/luarocks/manif.lua
usr/share/tarantool/luarocks/manif/
usr/share/tarantool/luarocks/manif/writer.lua
usr/share/tarantool/luarocks/pack.lua
usr/share/tarantool/luarocks/path.lua
usr/share/tarantool/luarocks/persist.lua
usr/share/tarantool/luarocks/queries.lua
usr/share/tarantool/luarocks/remove.lua
usr/share/tarantool/luarocks/repos.lua
usr/share/tarantool/luarocks/require.lua
usr/share/tarantool/luarocks/results.lua
usr/share/tarantool/luarocks/rockspecs.lua
usr/share/tarantool/luarocks/search.lua
usr/share/tarantool/luarocks/signing.lua
usr/share/tarantool/luarocks/test.lua
usr/share/tarantool/luarocks/test/
usr/share/tarantool/luarocks/test/busted.lua
usr/share/tarantool/luarocks/test/command.lua
usr/share/tarantool/luarocks/tools/
usr/share/tarantool/luarocks/tools/patch.lua
usr/share/tarantool/luarocks/tools/tar.lua
usr/share/tarantool/luarocks/tools/zip.lua
usr/share/tarantool/luarocks/type/
usr/share/tarantool/luarocks/type/manifest.lua
usr/share/tarantool/luarocks/type/rockspec.lua
usr/share/tarantool/luarocks/type_check.lua
usr/share/tarantool/luarocks/upload/
usr/share/tarantool/luarocks/upload/api.lua
usr/share/tarantool/luarocks/upload/multipart.lua
usr/share/tarantool/luarocks/util.lua
var/
var/lib/
var/lib/tarantool/
var/log/
var/log/tarantool/
