# extracted/mumble-server-1.5.735-1/files
%FILES%
etc/
etc/mumble/
etc/mumble/mumble-server.ini
usr/
usr/bin/
usr/bin/mumble-server
usr/bin/mumble-server-user-wrapper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mumble-server.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mumble-server.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mumble-server.conf
usr/share/
usr/share/licenses/
usr/share/licenses/mumble-server/
usr/share/licenses/mumble-server/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mumble-server-user-wrapper.1.gz
usr/share/man/man1/mumble-server.1.gz
usr/share/mumble-server/
usr/share/mumble-server/MumbleServer.ice
