# extracted/sshguard-2.4.3-2/files
%FILES%
etc/
etc/sshguard.conf
usr/
usr/bin/
usr/bin/sshguard
usr/lib/
usr/lib/sshguard/
usr/lib/sshguard/sshg-blocker
usr/lib/sshguard/sshg-fw-firewalld
usr/lib/sshguard/sshg-fw-hosts
usr/lib/sshguard/sshg-fw-ipfilter
usr/lib/sshguard/sshg-fw-ipfw
usr/lib/sshguard/sshg-fw-ipset
usr/lib/sshguard/sshg-fw-iptables
usr/lib/sshguard/sshg-fw-nft-sets
usr/lib/sshguard/sshg-fw-null
usr/lib/sshguard/sshg-fw-pf
usr/lib/sshguard/sshg-logtail
usr/lib/sshguard/sshg-parser
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sshguard.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sshguard.conf
usr/share/
usr/share/doc/
usr/share/doc/sshguard/
usr/share/doc/sshguard/net.sshguard.plist
usr/share/doc/sshguard/sshguard.conf.sample
usr/share/doc/sshguard/sshguard.service
usr/share/doc/sshguard/whitelistfile.example
usr/share/man/
usr/share/man/man7/
usr/share/man/man7/sshguard-setup.7.gz
usr/share/man/man8/
usr/share/man/man8/sshguard.8.gz
