# extracted/loki-3.2.1-1/files
%FILES%
etc/
etc/loki/
etc/loki/loki.yaml
usr/
usr/bin/
usr/bin/loki
usr/bin/lokitool
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/loki.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/loki.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/loki.conf
