# extracted/lldpd-1.0.18-2/files
%FILES%
etc/
etc/lldpd.conf
etc/lldpd.d/
etc/lldpd.d/README.conf
usr/
usr/bin/
usr/bin/lldpcli
usr/bin/lldpctl
usr/bin/lldpd
usr/include/
usr/include/lldp-const.h
usr/include/lldpctl.h
usr/lib/
usr/lib/liblldpctl.so
usr/lib/liblldpctl.so.4
usr/lib/liblldpctl.so.4.9.1
usr/lib/pkgconfig/
usr/lib/pkgconfig/lldpctl.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lldpd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/lldpd.sysusers.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/lldpd.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/lldpcli
usr/share/doc/
usr/share/doc/lldpd/
usr/share/doc/lldpd/CONTRIBUTE.md
usr/share/doc/lldpd/ChangeLog
usr/share/doc/lldpd/LICENSE
usr/share/doc/lldpd/NEWS
usr/share/doc/lldpd/README.md
usr/share/licenses/
usr/share/licenses/lldpd/
usr/share/licenses/lldpd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/lldpcli.8.gz
usr/share/man/man8/lldpctl.8.gz
usr/share/man/man8/lldpd.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_lldpcli
