# extracted/mpdscribble-0.25-2/files
%FILES%
usr/
usr/bin/
usr/bin/mpdscribble
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mpdscribble.service
usr/lib/systemd/system/mpdscribble@.service
usr/lib/systemd/user/
usr/lib/systemd/user/mpdscribble.service
usr/share/
usr/share/doc/
usr/share/doc/mpdscribble/
usr/share/doc/mpdscribble/AUTHORS
usr/share/doc/mpdscribble/COPYING
usr/share/doc/mpdscribble/NEWS
usr/share/doc/mpdscribble/README.rst
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mpdscribble.1.gz
usr/share/mpdscribble/
usr/share/mpdscribble/mpdscribble.conf.example
var/
var/cache/
var/cache/mpdscribble/
var/cache/mpdscribble/mpdscribble.cache
