# extracted/interception-tools-0.6.8-8/files
%FILES%
etc/
etc/interception/
etc/interception/udevmon.d/
usr/
usr/bin/
usr/bin/intercept
usr/bin/mux
usr/bin/udevmon
usr/bin/uinput
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/udevmon.service
usr/share/
usr/share/doc/
usr/share/doc/interception-tools/
usr/share/doc/interception-tools/README.md
