# extracted/dante-1.4.3-4/files
%FILES%
etc/
etc/sockd.conf
etc/socks.conf
usr/
usr/bin/
usr/bin/sockd
usr/bin/socksify
usr/include/
usr/include/socks.h
usr/lib/
usr/lib/libdsocks.so
usr/lib/libsocks.so
usr/lib/libsocks.so.0
usr/lib/libsocks.so.0.1.1
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sockd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sockd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/dante/
usr/share/licenses/dante/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/socksify.1.gz
usr/share/man/man5/
usr/share/man/man5/sockd.conf.5.gz
usr/share/man/man5/socks.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/sockd.8.gz
