# extracted/opendkim-2.11.0beta-8/files
%FILES%
usr/
usr/bin/
usr/bin/opendkim
usr/bin/opendkim-genkey
usr/bin/opendkim-genzone
usr/bin/opendkim-testkey
usr/bin/opendkim-testmsg
usr/include/
usr/include/opendkim/
usr/include/opendkim/dkim.h
usr/lib/
usr/lib/libopendkim.so
usr/lib/libopendkim.so.11
usr/lib/libopendkim.so.11.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/opendkim.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opendkim.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/opendkim.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/opendkim.conf
usr/share/
usr/share/doc/
usr/share/doc/opendkim/
usr/share/doc/opendkim/FEATURES
usr/share/doc/opendkim/KNOWNBUGS
usr/share/doc/opendkim/LICENSE
usr/share/doc/opendkim/LICENSE.Sendmail
usr/share/doc/opendkim/README
usr/share/doc/opendkim/README.LDAP
usr/share/doc/opendkim/README.SQL
usr/share/doc/opendkim/README.specs.html
usr/share/doc/opendkim/RELEASE_NOTES.Sendmail
usr/share/doc/opendkim/authheaders-check-setup-hook.lua
usr/share/doc/opendkim/chroot
usr/share/doc/opendkim/convert_keylist.sh
usr/share/doc/opendkim/dkim.html
usr/share/doc/opendkim/dkim_add_querymethod.html
usr/share/doc/opendkim/dkim_add_xtag.html
usr/share/doc/opendkim/dkim_alg_t.html
usr/share/doc/opendkim/dkim_atps_check.html
usr/share/doc/opendkim/dkim_atps_t.html
usr/share/doc/opendkim/dkim_body.html
usr/share/doc/opendkim/dkim_canon_t.html
usr/share/doc/opendkim/dkim_cbstat.html
usr/share/doc/opendkim/dkim_chunk.html
usr/share/doc/opendkim/dkim_close.html
usr/share/doc/opendkim/dkim_diffheaders.html
usr/share/doc/opendkim/dkim_dns_close.html
usr/share/doc/opendkim/dkim_dns_config.html
usr/share/doc/opendkim/dkim_dns_init.html
usr/share/doc/opendkim/dkim_dns_nslist.html
usr/share/doc/opendkim/dkim_dns_set_close.html
usr/share/doc/opendkim/dkim_dns_set_config.html
usr/share/doc/opendkim/dkim_dns_set_init.html
usr/share/doc/opendkim/dkim_dns_set_nslist.html
usr/share/doc/opendkim/dkim_dns_set_query_cancel.html
usr/share/doc/opendkim/dkim_dns_set_query_service.html
usr/share/doc/opendkim/dkim_dns_set_query_start.html
usr/share/doc/opendkim/dkim_dns_set_query_waitreply.html
usr/share/doc/opendkim/dkim_dns_set_trustanchor.html
usr/share/doc/opendkim/dkim_dns_trustanchor.html
usr/share/doc/opendkim/dkim_dnssec.html
usr/share/doc/opendkim/dkim_eoh.html
usr/share/doc/opendkim/dkim_eom.html
usr/share/doc/opendkim/dkim_flush_cache.html
usr/share/doc/opendkim/dkim_free.html
usr/share/doc/opendkim/dkim_get_msgdate.html
usr/share/doc/opendkim/dkim_get_reputation.html
usr/share/doc/opendkim/dkim_get_signer.html
usr/share/doc/opendkim/dkim_get_sigsubstring.html
usr/share/doc/opendkim/dkim_get_user_context.html
usr/share/doc/opendkim/dkim_getcachestats.html
usr/share/doc/opendkim/dkim_getdomain.html
usr/share/doc/opendkim/dkim_geterror.html
usr/share/doc/opendkim/dkim_getid.html
usr/share/doc/opendkim/dkim_getmode.html
usr/share/doc/opendkim/dkim_getpartial.html
usr/share/doc/opendkim/dkim_getresultstr.html
usr/share/doc/opendkim/dkim_getsighdr.html
usr/share/doc/opendkim/dkim_getsighdr_d.html
usr/share/doc/opendkim/dkim_getsiglist.html
usr/share/doc/opendkim/dkim_getsignature.html
usr/share/doc/opendkim/dkim_getsslbuf.html
usr/share/doc/opendkim/dkim_getuser.html
usr/share/doc/opendkim/dkim_header.html
usr/share/doc/opendkim/dkim_init.html
usr/share/doc/opendkim/dkim_key_syntax.html
usr/share/doc/opendkim/dkim_lib.html
usr/share/doc/opendkim/dkim_libfeature.html
usr/share/doc/opendkim/dkim_libversion.html
usr/share/doc/opendkim/dkim_mail_parse.html
usr/share/doc/opendkim/dkim_minbody.html
usr/share/doc/opendkim/dkim_ohdrs.html
usr/share/doc/opendkim/dkim_options.html
usr/share/doc/opendkim/dkim_param_t.html
usr/share/doc/opendkim/dkim_privkey_load.html
usr/share/doc/opendkim/dkim_qi_getname.html
usr/share/doc/opendkim/dkim_qi_gettype.html
usr/share/doc/opendkim/dkim_query_t.html
usr/share/doc/opendkim/dkim_queryinfo.html
usr/share/doc/opendkim/dkim_resign.html
usr/share/doc/opendkim/dkim_set_dns_callback.html
usr/share/doc/opendkim/dkim_set_final.html
usr/share/doc/opendkim/dkim_set_key_lookup.html
usr/share/doc/opendkim/dkim_set_margin.html
usr/share/doc/opendkim/dkim_set_prescreen.html
usr/share/doc/opendkim/dkim_set_signature_handle.html
usr/share/doc/opendkim/dkim_set_signature_handle_free.html
usr/share/doc/opendkim/dkim_set_signature_tagvalues.html
usr/share/doc/opendkim/dkim_set_signer.html
usr/share/doc/opendkim/dkim_set_trust_anchor.html
usr/share/doc/opendkim/dkim_set_user_context.html
usr/share/doc/opendkim/dkim_setpartial.html
usr/share/doc/opendkim/dkim_sig_getbh.html
usr/share/doc/opendkim/dkim_sig_getcanonlen.html
usr/share/doc/opendkim/dkim_sig_getcanons.html
usr/share/doc/opendkim/dkim_sig_getcontext.html
usr/share/doc/opendkim/dkim_sig_getdnssec.html
usr/share/doc/opendkim/dkim_sig_getdomain.html
usr/share/doc/opendkim/dkim_sig_geterror.html
usr/share/doc/opendkim/dkim_sig_geterrorstr.html
usr/share/doc/opendkim/dkim_sig_getflags.html
usr/share/doc/opendkim/dkim_sig_gethashes.html
usr/share/doc/opendkim/dkim_sig_getidentity.html
usr/share/doc/opendkim/dkim_sig_getkeysize.html
usr/share/doc/opendkim/dkim_sig_getqueries.html
usr/share/doc/opendkim/dkim_sig_getreportinfo.html
usr/share/doc/opendkim/dkim_sig_getselector.html
usr/share/doc/opendkim/dkim_sig_getsignalg.html
usr/share/doc/opendkim/dkim_sig_getsignedhdrs.html
usr/share/doc/opendkim/dkim_sig_getsigntime.html
usr/share/doc/opendkim/dkim_sig_getsslbuf.html
usr/share/doc/opendkim/dkim_sig_gettagvalue.html
usr/share/doc/opendkim/dkim_sig_hdrsigned.html
usr/share/doc/opendkim/dkim_sig_ignore.html
usr/share/doc/opendkim/dkim_sig_process.html
usr/share/doc/opendkim/dkim_sig_seterror.html
usr/share/doc/opendkim/dkim_sig_syntax.html
usr/share/doc/opendkim/dkim_sigerror.html
usr/share/doc/opendkim/dkim_siginfo.html
usr/share/doc/opendkim/dkim_sigkey_t.html
usr/share/doc/opendkim/dkim_sign.html
usr/share/doc/opendkim/dkim_signhdrs.html
usr/share/doc/opendkim/dkim_ssl_version.html
usr/share/doc/opendkim/dkim_stat.html
usr/share/doc/opendkim/dkim_verify.html
usr/share/doc/opendkim/dns.html
usr/share/doc/opendkim/example.com.ldif
usr/share/doc/opendkim/index.html
usr/share/doc/opendkim/opendkim
usr/share/doc/opendkim/opendkim-default-keygen
usr/share/doc/opendkim/opendkim-fixipaddrs.pl
usr/share/doc/opendkim/opendkim-spam-ignore-cmdswitch.patch
usr/share/doc/opendkim/opendkim.conf.sample
usr/share/doc/opendkim/opendkim.conf.simple
usr/share/doc/opendkim/opendkim.conf.simple-verify
usr/share/doc/opendkim/opendkim.ldif
usr/share/doc/opendkim/opendkim.schema
usr/share/doc/opendkim/opendkim.service
usr/share/doc/opendkim/opendkim.spec.in
usr/share/doc/opendkim/opendkim.xml
usr/share/doc/opendkim/overview.html
usr/share/doc/opendkim/repute.py
usr/share/doc/opendkim/stats.lua
usr/share/licenses/
usr/share/licenses/opendkim/
usr/share/licenses/opendkim/LICENSE
usr/share/licenses/opendkim/LICENSE.Sendmail
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/opendkim.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/opendkim-genkey.8.gz
usr/share/man/man8/opendkim-genzone.8.gz
usr/share/man/man8/opendkim-testkey.8.gz
usr/share/man/man8/opendkim-testmsg.8.gz
usr/share/man/man8/opendkim.8.gz
