# extracted/incron-0.5.12-6/files
%FILES%
etc/
etc/incron.conf
etc/incron.d/
usr/
usr/bin/
usr/bin/incrond
usr/bin/incrontab
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/incrond.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/incrond.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/incrond.conf
usr/share/
usr/share/doc/
usr/share/doc/incron-0.5.12/
usr/share/doc/incron-0.5.12/incron.conf.example
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/incrontab.1.gz
usr/share/man/man5/
usr/share/man/man5/incron.conf.5.gz
usr/share/man/man5/incrontab.5.gz
usr/share/man/man8/
usr/share/man/man8/incrond.8.gz
