# extracted/wsdd-0.8-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/wsdd
usr/
usr/bin/
usr/bin/wsdd
usr/lib/
usr/lib/firewalld/
usr/lib/firewalld/services/
usr/lib/firewalld/services/wsdd-http.xml
usr/lib/firewalld/services/wsdd.xml
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wsdd.service
usr/share/
usr/share/doc/
usr/share/doc/wsdd/
usr/share/doc/wsdd/AUTHORS
usr/share/doc/wsdd/README.md
usr/share/licenses/
usr/share/licenses/wsdd/
usr/share/licenses/wsdd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/wsdd.8.gz
