# extracted/iodine-0.8.0-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/iodined
usr/
usr/bin/
usr/bin/iodine
usr/bin/iodined
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/iodined.service
usr/lib/systemd/system/iodined.socket
usr/share/
usr/share/doc/
usr/share/doc/iodine/
usr/share/doc/iodine/README.md
usr/share/licenses/
usr/share/licenses/iodine/
usr/share/licenses/iodine/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/iodine.8.gz
