# extracted/taskd-1.1.0-6/files
%FILES%
usr/
usr/bin/
usr/bin/taskd
usr/bin/taskdctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/taskd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/taskd.conf
usr/share/
usr/share/doc/
usr/share/doc/taskd/
usr/share/doc/taskd/AUTHORS
usr/share/doc/taskd/COPYING
usr/share/doc/taskd/ChangeLog
usr/share/doc/taskd/INSTALL
usr/share/doc/taskd/NEWS
usr/share/doc/taskd/pki/
usr/share/doc/taskd/pki/README
usr/share/doc/taskd/pki/generate
usr/share/doc/taskd/pki/generate.ca
usr/share/doc/taskd/pki/generate.client
usr/share/doc/taskd/pki/generate.crl
usr/share/doc/taskd/pki/generate.server
usr/share/doc/taskd/pki/vars
usr/share/doc/taskd/scripts/
usr/share/doc/taskd/scripts/l10n
usr/share/doc/taskd/scripts/profile.py
usr/share/doc/taskd/scripts/setup_server.bash
usr/share/doc/taskd/scripts/systemd/
usr/share/doc/taskd/scripts/systemd/taskd.service
usr/share/licenses/
usr/share/licenses/taskd/
usr/share/licenses/taskd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/taskd.1.gz
usr/share/man/man1/taskdctl.1.gz
usr/share/man/man5/
usr/share/man/man5/taskdrc.5.gz
