# extracted/bitcoin-daemon-28.0-1/files
%FILES%
etc/
etc/bitcoin/
etc/bitcoin/bitcoin.conf
usr/
usr/bin/
usr/bin/bitcoin-cli
usr/bin/bitcoind
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bitcoind.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/bitcoin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/bitcoin.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/bitcoin-cli
usr/share/bash-completion/completions/bitcoind
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/bitcoin-cli.fish
usr/share/fish/vendor_completions.d/bitcoind.fish
usr/share/licenses/
usr/share/licenses/bitcoin-daemon/
usr/share/licenses/bitcoin-daemon/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/bitcoin-cli.1.gz
usr/share/man/man1/bitcoind.1.gz
