# extracted/sauerbraten-2020_12_29-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/sauerbraten
opt/
opt/sauerbraten/
opt/sauerbraten/bin/
opt/sauerbraten/bin/linux_client
opt/sauerbraten/bin/linux_server
usr/
usr/bin/
usr/bin/sauerbraten-client
usr/bin/sauerbraten-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sauerbraten.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/sauerbraten.conf
usr/share/
usr/share/applications/
usr/share/applications/sauerbraten.desktop
usr/share/licenses/
usr/share/licenses/sauerbraten/
usr/share/licenses/sauerbraten/license.txt
usr/share/pixmaps/
usr/share/pixmaps/sauerbraten.png
