# extracted/iperf-2.2.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/iperf
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/iperf-tcp.service
usr/lib/systemd/system/iperf-udp.service
usr/share/
usr/share/licenses/
usr/share/licenses/iperf/
usr/share/licenses/iperf/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/iperf.1.gz
