# extracted/nix-2.24.7-2/files
%FILES%
etc/
etc/nix/
etc/nix/nix.conf
etc/profile.d/
etc/profile.d/nix-daemon.fish
etc/profile.d/nix-daemon.sh
usr/
usr/bin/
usr/bin/nix
usr/bin/nix-build
usr/bin/nix-channel
usr/bin/nix-collect-garbage
usr/bin/nix-copy-closure
usr/bin/nix-daemon
usr/bin/nix-env
usr/bin/nix-hash
usr/bin/nix-instantiate
usr/bin/nix-prefetch-url
usr/bin/nix-shell
usr/bin/nix-store
usr/include/
usr/include/nix/
usr/include/nix/abstract-setting-to-json.hh
usr/include/nix/ansicolor.hh
usr/include/nix/archive.hh
usr/include/nix/args.hh
usr/include/nix/args/
usr/include/nix/args/root.hh
usr/include/nix/attr-path.hh
usr/include/nix/attr-set.hh
usr/include/nix/attrs.hh
usr/include/nix/binary-cache-store.hh
usr/include/nix/build-result.hh
usr/include/nix/build/
usr/include/nix/build/derivation-goal.hh
usr/include/nix/build/drv-output-substitution-goal.hh
usr/include/nix/build/goal.hh
usr/include/nix/build/substitution-goal.hh
usr/include/nix/build/worker.hh
usr/include/nix/built-path.hh
usr/include/nix/builtins.hh
usr/include/nix/builtins/
usr/include/nix/builtins/buildenv.hh
usr/include/nix/ca-specific-schema.sql.gen.hh
usr/include/nix/cache.hh
usr/include/nix/call-flake.nix.gen.hh
usr/include/nix/callback.hh
usr/include/nix/canon-path.hh
usr/include/nix/chunked-vector.hh
usr/include/nix/closure.hh
usr/include/nix/command-installable-value.hh
usr/include/nix/command.hh
usr/include/nix/common-args.hh
usr/include/nix/common-eval-args.hh
usr/include/nix/common-protocol-impl.hh
usr/include/nix/common-protocol.hh
usr/include/nix/common-ssh-store-config.hh
usr/include/nix/comparator.hh
usr/include/nix/compatibility-settings.hh
usr/include/nix/compression.hh
usr/include/nix/compute-levels.hh
usr/include/nix/config-global.hh
usr/include/nix/config-impl.hh
usr/include/nix/config.h
usr/include/nix/config.hh
usr/include/nix/content-address.hh
usr/include/nix/current-process.hh
usr/include/nix/daemon.hh
usr/include/nix/derivations.hh
usr/include/nix/derived-path-map.hh
usr/include/nix/derived-path.hh
usr/include/nix/downstream-placeholder.hh
usr/include/nix/editor-for.hh
usr/include/nix/english.hh
usr/include/nix/environment-variables.hh
usr/include/nix/error.hh
usr/include/nix/eval-cache.hh
usr/include/nix/eval-error.hh
usr/include/nix/eval-gc.hh
usr/include/nix/eval-inline.hh
usr/include/nix/eval-settings.hh
usr/include/nix/eval.hh
usr/include/nix/exit.hh
usr/include/nix/experimental-features.hh
usr/include/nix/fetch-settings.hh
usr/include/nix/fetch-to-store.hh
usr/include/nix/fetchers.hh
usr/include/nix/fetchurl.nix.gen.hh
usr/include/nix/file-content-address.hh
usr/include/nix/file-descriptor.hh
usr/include/nix/file-path-impl.hh
usr/include/nix/file-path.hh
usr/include/nix/file-system.hh
usr/include/nix/filetransfer.hh
usr/include/nix/filtering-source-accessor.hh
usr/include/nix/finally.hh
usr/include/nix/flake/
usr/include/nix/flake/flake.hh
usr/include/nix/flake/flakeref.hh
usr/include/nix/flake/lockfile.hh
usr/include/nix/flake/settings.hh
usr/include/nix/flake/url-name.hh
usr/include/nix/fmt.hh
usr/include/nix/fs-sink.hh
usr/include/nix/function-trace.hh
usr/include/nix/gc-small-vector.hh
usr/include/nix/gc-store.hh
usr/include/nix/get-drvs.hh
usr/include/nix/git-utils.hh
usr/include/nix/git.hh
usr/include/nix/globals.hh
usr/include/nix/hash.hh
usr/include/nix/hilite.hh
usr/include/nix/http-binary-cache-store.hh
usr/include/nix/imported-drv-to-derivation.nix.gen.hh
usr/include/nix/indirect-root-store.hh
usr/include/nix/installable-attr-path.hh
usr/include/nix/installable-derived-path.hh
usr/include/nix/installable-flake.hh
usr/include/nix/installable-value.hh
usr/include/nix/installables.hh
usr/include/nix/json-impls.hh
usr/include/nix/json-to-value.hh
usr/include/nix/json-utils.hh
usr/include/nix/keys.hh
usr/include/nix/legacy-ssh-store.hh
usr/include/nix/legacy.hh
usr/include/nix/length-prefixed-protocol-helper.hh
usr/include/nix/lexer-helpers.hh
usr/include/nix/lexer-tab.hh
usr/include/nix/local-binary-cache-store.hh
usr/include/nix/local-fs-store.hh
usr/include/nix/local-overlay-store.hh
usr/include/nix/local-store.hh
usr/include/nix/log-store.hh
usr/include/nix/loggers.hh
usr/include/nix/logging.hh
usr/include/nix/lru-cache.hh
usr/include/nix/machines.hh
usr/include/nix/make-content-addressed.hh
usr/include/nix/markdown.hh
usr/include/nix/memory-source-accessor.hh
usr/include/nix/misc-store-flags.hh
usr/include/nix/monitor-fd.hh
usr/include/nix/mounted-source-accessor.hh
usr/include/nix/muxable-pipe.hh
usr/include/nix/names.hh
usr/include/nix/nar-accessor.hh
usr/include/nix/nar-info-disk-cache.hh
usr/include/nix/nar-info.hh
usr/include/nix/network-proxy.hh
usr/include/nix/nix_api_expr.h
usr/include/nix/nix_api_external.h
usr/include/nix/nix_api_main.h
usr/include/nix/nix_api_store.h
usr/include/nix/nix_api_util.h
usr/include/nix/nix_api_value.h
usr/include/nix/nixexpr.hh
usr/include/nix/outputs-spec.hh
usr/include/nix/parsed-derivations.hh
usr/include/nix/parser-state.hh
usr/include/nix/parser-tab.hh
usr/include/nix/path-info.hh
usr/include/nix/path-references.hh
usr/include/nix/path-regex.hh
usr/include/nix/path-with-outputs.hh
usr/include/nix/path.hh
usr/include/nix/pathlocks.hh
usr/include/nix/plugin.hh
usr/include/nix/pool.hh
usr/include/nix/pos-idx.hh
usr/include/nix/pos-table.hh
usr/include/nix/position.hh
usr/include/nix/posix-fs-canonicalise.hh
usr/include/nix/posix-source-accessor.hh
usr/include/nix/primops.hh
usr/include/nix/print-ambiguous.hh
usr/include/nix/print-options.hh
usr/include/nix/print.hh
usr/include/nix/processes.hh
usr/include/nix/profiles.hh
usr/include/nix/progress-bar.hh
usr/include/nix/realisation.hh
usr/include/nix/ref.hh
usr/include/nix/references.hh
usr/include/nix/regex-combinators.hh
usr/include/nix/registry.hh
usr/include/nix/remote-fs-accessor.hh
usr/include/nix/remote-store-connection.hh
usr/include/nix/remote-store.hh
usr/include/nix/repair-flag.hh
usr/include/nix/repl-exit-status.hh
usr/include/nix/repl-interacter.hh
usr/include/nix/repl.hh
usr/include/nix/s3-binary-cache-store.hh
usr/include/nix/s3.hh
usr/include/nix/schema.sql.gen.hh
usr/include/nix/search-path.hh
usr/include/nix/serialise.hh
usr/include/nix/serve-protocol-connection.hh
usr/include/nix/serve-protocol-impl.hh
usr/include/nix/serve-protocol.hh
usr/include/nix/shared.hh
usr/include/nix/signals-impl.hh
usr/include/nix/signals.hh
usr/include/nix/signature/
usr/include/nix/signature/local-keys.hh
usr/include/nix/signature/signer.hh
usr/include/nix/source-accessor.hh
usr/include/nix/source-path.hh
usr/include/nix/split.hh
usr/include/nix/sqlite.hh
usr/include/nix/ssh-store.hh
usr/include/nix/ssh.hh
usr/include/nix/std-hash.hh
usr/include/nix/store-api.hh
usr/include/nix/store-cast.hh
usr/include/nix/store-dir-config.hh
usr/include/nix/store-path-accessor.hh
usr/include/nix/store-reference.hh
usr/include/nix/strings-inline.hh
usr/include/nix/strings.hh
usr/include/nix/suggestions.hh
usr/include/nix/symbol-table.hh
usr/include/nix/sync.hh
usr/include/nix/tarball.hh
usr/include/nix/tarfile.hh
usr/include/nix/terminal.hh
usr/include/nix/thread-pool.hh
usr/include/nix/topo-sort.hh
usr/include/nix/types.hh
usr/include/nix/uds-remote-store.hh
usr/include/nix/unix-domain-socket.hh
usr/include/nix/url-parts.hh
usr/include/nix/url.hh
usr/include/nix/user-lock.hh
usr/include/nix/users.hh
usr/include/nix/util.hh
usr/include/nix/value-to-json.hh
usr/include/nix/value-to-xml.hh
usr/include/nix/value.hh
usr/include/nix/value/
usr/include/nix/value/context.hh
usr/include/nix/variant-wrapper.hh
usr/include/nix/worker-protocol-connection.hh
usr/include/nix/worker-protocol-impl.hh
usr/include/nix/worker-protocol.hh
usr/include/nix/xml-writer.hh
usr/lib/
usr/lib/libnixcmd.so
usr/lib/libnixexpr.so
usr/lib/libnixexprc.so
usr/lib/libnixfetchers.so
usr/lib/libnixflake.so
usr/lib/libnixmain.so
usr/lib/libnixstore.so
usr/lib/libnixstorec.so
usr/lib/libnixutil.so
usr/lib/libnixutilc.so
usr/lib/nix/
usr/lib/nix/nix/
usr/lib/nix/nix/build-remote
usr/lib/pkgconfig/
usr/lib/pkgconfig/nix-cmd.pc
usr/lib/pkgconfig/nix-expr-c.pc
usr/lib/pkgconfig/nix-expr.pc
usr/lib/pkgconfig/nix-flake.pc
usr/lib/pkgconfig/nix-main.pc
usr/lib/pkgconfig/nix-store-c.pc
usr/lib/pkgconfig/nix-store.pc
usr/lib/pkgconfig/nix-util.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nix-daemon.service
usr/lib/systemd/system/nix-daemon.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/nix-daemon.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nix-daemon.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/nix
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/nix.fish
usr/share/licenses/
usr/share/licenses/nix/
usr/share/licenses/nix/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/nix-build.1.gz
usr/share/man/man1/nix-channel.1.gz
usr/share/man/man1/nix-collect-garbage.1.gz
usr/share/man/man1/nix-copy-closure.1.gz
usr/share/man/man1/nix-env-delete-generations.1.gz
usr/share/man/man1/nix-env-install.1.gz
usr/share/man/man1/nix-env-list-generations.1.gz
usr/share/man/man1/nix-env-query.1.gz
usr/share/man/man1/nix-env-rollback.1.gz
usr/share/man/man1/nix-env-set-flag.1.gz
usr/share/man/man1/nix-env-set.1.gz
usr/share/man/man1/nix-env-switch-generation.1.gz
usr/share/man/man1/nix-env-switch-profile.1.gz
usr/share/man/man1/nix-env-uninstall.1.gz
usr/share/man/man1/nix-env-upgrade.1.gz
usr/share/man/man1/nix-env.1.gz
usr/share/man/man1/nix-hash.1.gz
usr/share/man/man1/nix-instantiate.1.gz
usr/share/man/man1/nix-prefetch-url.1.gz
usr/share/man/man1/nix-shell.1.gz
usr/share/man/man1/nix-store-add-fixed.1.gz
usr/share/man/man1/nix-store-add.1.gz
usr/share/man/man1/nix-store-delete.1.gz
usr/share/man/man1/nix-store-dump-db.1.gz
usr/share/man/man1/nix-store-dump.1.gz
usr/share/man/man1/nix-store-export.1.gz
usr/share/man/man1/nix-store-gc.1.gz
usr/share/man/man1/nix-store-generate-binary-cache-key.1.gz
usr/share/man/man1/nix-store-import.1.gz
usr/share/man/man1/nix-store-load-db.1.gz
usr/share/man/man1/nix-store-optimise.1.gz
usr/share/man/man1/nix-store-print-env.1.gz
usr/share/man/man1/nix-store-query.1.gz
usr/share/man/man1/nix-store-read-log.1.gz
usr/share/man/man1/nix-store-realise.1.gz
usr/share/man/man1/nix-store-repair-path.1.gz
usr/share/man/man1/nix-store-restore.1.gz
usr/share/man/man1/nix-store-serve.1.gz
usr/share/man/man1/nix-store-verify-path.1.gz
usr/share/man/man1/nix-store-verify.1.gz
usr/share/man/man1/nix-store.1.gz
usr/share/man/man1/nix.1.gz
usr/share/man/man1/nix3-build.1.gz
usr/share/man/man1/nix3-bundle.1.gz
usr/share/man/man1/nix3-config-check.1.gz
usr/share/man/man1/nix3-config-show.1.gz
usr/share/man/man1/nix3-config.1.gz
usr/share/man/man1/nix3-copy.1.gz
usr/share/man/man1/nix3-daemon.1.gz
usr/share/man/man1/nix3-derivation-add.1.gz
usr/share/man/man1/nix3-derivation-show.1.gz
usr/share/man/man1/nix3-derivation.1.gz
usr/share/man/man1/nix3-develop.1.gz
usr/share/man/man1/nix3-edit.1.gz
usr/share/man/man1/nix3-env-shell.1.gz
usr/share/man/man1/nix3-env.1.gz
usr/share/man/man1/nix3-eval.1.gz
usr/share/man/man1/nix3-flake-archive.1.gz
usr/share/man/man1/nix3-flake-check.1.gz
usr/share/man/man1/nix3-flake-clone.1.gz
usr/share/man/man1/nix3-flake-info.1.gz
usr/share/man/man1/nix3-flake-init.1.gz
usr/share/man/man1/nix3-flake-lock.1.gz
usr/share/man/man1/nix3-flake-metadata.1.gz
usr/share/man/man1/nix3-flake-new.1.gz
usr/share/man/man1/nix3-flake-prefetch.1.gz
usr/share/man/man1/nix3-flake-show.1.gz
usr/share/man/man1/nix3-flake-update.1.gz
usr/share/man/man1/nix3-flake.1.gz
usr/share/man/man1/nix3-fmt.1.gz
usr/share/man/man1/nix3-hash-convert.1.gz
usr/share/man/man1/nix3-hash-file.1.gz
usr/share/man/man1/nix3-hash-path.1.gz
usr/share/man/man1/nix3-hash-to-base16.1.gz
usr/share/man/man1/nix3-hash-to-base32.1.gz
usr/share/man/man1/nix3-hash-to-base64.1.gz
usr/share/man/man1/nix3-hash-to-sri.1.gz
usr/share/man/man1/nix3-hash.1.gz
usr/share/man/man1/nix3-help-stores.1.gz
usr/share/man/man1/nix3-help.1.gz
usr/share/man/man1/nix3-key-convert-secret-to-public.1.gz
usr/share/man/man1/nix3-key-generate-secret.1.gz
usr/share/man/man1/nix3-key.1.gz
usr/share/man/man1/nix3-log.1.gz
usr/share/man/man1/nix3-manpages.gz
usr/share/man/man1/nix3-nar-cat.1.gz
usr/share/man/man1/nix3-nar-dump-path.1.gz
usr/share/man/man1/nix3-nar-ls.1.gz
usr/share/man/man1/nix3-nar-pack.1.gz
usr/share/man/man1/nix3-nar.1.gz
usr/share/man/man1/nix3-path-info.1.gz
usr/share/man/man1/nix3-print-dev-env.1.gz
usr/share/man/man1/nix3-profile-diff-closures.1.gz
usr/share/man/man1/nix3-profile-history.1.gz
usr/share/man/man1/nix3-profile-install.1.gz
usr/share/man/man1/nix3-profile-list.1.gz
usr/share/man/man1/nix3-profile-remove.1.gz
usr/share/man/man1/nix3-profile-rollback.1.gz
usr/share/man/man1/nix3-profile-upgrade.1.gz
usr/share/man/man1/nix3-profile-wipe-history.1.gz
usr/share/man/man1/nix3-profile.1.gz
usr/share/man/man1/nix3-realisation-info.1.gz
usr/share/man/man1/nix3-realisation.1.gz
usr/share/man/man1/nix3-registry-add.1.gz
usr/share/man/man1/nix3-registry-list.1.gz
usr/share/man/man1/nix3-registry-pin.1.gz
usr/share/man/man1/nix3-registry-remove.1.gz
usr/share/man/man1/nix3-registry.1.gz
usr/share/man/man1/nix3-repl.1.gz
usr/share/man/man1/nix3-run.1.gz
usr/share/man/man1/nix3-search.1.gz
usr/share/man/man1/nix3-store-add-file.1.gz
usr/share/man/man1/nix3-store-add-path.1.gz
usr/share/man/man1/nix3-store-add.1.gz
usr/share/man/man1/nix3-store-cat.1.gz
usr/share/man/man1/nix3-store-copy-log.1.gz
usr/share/man/man1/nix3-store-copy-sigs.1.gz
usr/share/man/man1/nix3-store-delete.1.gz
usr/share/man/man1/nix3-store-diff-closures.1.gz
usr/share/man/man1/nix3-store-dump-path.1.gz
usr/share/man/man1/nix3-store-gc.1.gz
usr/share/man/man1/nix3-store-info.1.gz
usr/share/man/man1/nix3-store-ls.1.gz
usr/share/man/man1/nix3-store-make-content-addressed.1.gz
usr/share/man/man1/nix3-store-optimise.1.gz
usr/share/man/man1/nix3-store-path-from-hash-part.1.gz
usr/share/man/man1/nix3-store-ping.1.gz
usr/share/man/man1/nix3-store-prefetch-file.1.gz
usr/share/man/man1/nix3-store-repair.1.gz
usr/share/man/man1/nix3-store-sign.1.gz
usr/share/man/man1/nix3-store-verify.1.gz
usr/share/man/man1/nix3-store.1.gz
usr/share/man/man1/nix3-upgrade-nix.1.gz
usr/share/man/man1/nix3-why-depends.1.gz
usr/share/man/man5/
usr/share/man/man5/nix-profiles.5.gz
usr/share/man/man5/nix.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/nix-daemon.8.gz
usr/share/user-tmpfiles.d/
usr/share/user-tmpfiles.d/nix-daemon.conf
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_nix
usr/share/zsh/site-functions/run-help-nix
