# extracted/linuxconsole-1.8.1-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/inputattach
usr/
usr/bin/
usr/bin/evdev-joystick
usr/bin/ffcfstress
usr/bin/ffmvforce
usr/bin/ffset
usr/bin/fftest
usr/bin/inputattach
usr/bin/inputattachctl
usr/bin/jscal
usr/bin/jscal-restore
usr/bin/jscal-store
usr/bin/jstest
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/inputattach.service
usr/lib/udev/
usr/lib/udev/js-set-enum-leds
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/80-stelladaptor-joystick.rules
usr/share/
usr/share/joystick/
usr/share/joystick/extract
usr/share/joystick/filter
usr/share/joystick/ident
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/evdev-joystick.1.gz
usr/share/man/man1/ffcfstress.1.gz
usr/share/man/man1/ffmvforce.1.gz
usr/share/man/man1/ffset.1.gz
usr/share/man/man1/fftest.1.gz
usr/share/man/man1/inputattach.1.gz
usr/share/man/man1/jscal-restore.1.gz
usr/share/man/man1/jscal-store.1.gz
usr/share/man/man1/jscal.1.gz
usr/share/man/man1/jstest.1.gz
