# extracted/prometheus-smokeping-prober-0.8.1-1/files
%FILES%
etc/
etc/prometheus/
etc/prometheus/smokeping_prober.yml
usr/
usr/bin/
usr/bin/prometheus-smokeping-prober
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-smokeping-prober.service
usr/share/
usr/share/doc/
usr/share/doc/prometheus-smokeping-prober/
usr/share/doc/prometheus-smokeping-prober/CHANGELOG.md
usr/share/doc/prometheus-smokeping-prober/README.md
usr/share/doc/prometheus-smokeping-prober/dashboard.json
usr/share/doc/prometheus-smokeping-prober/example-rules.yml
usr/share/licenses/
usr/share/licenses/prometheus-smokeping-prober/
usr/share/licenses/prometheus-smokeping-prober/LICENSE
