# extracted/motion-4.7.0-2/files
%FILES%
etc/
etc/motion/
etc/motion/camera1-dist.conf
etc/motion/camera2-dist.conf
etc/motion/camera3-dist.conf
etc/motion/camera4-dist.conf
etc/motion/motion-dist.conf
etc/motion/motion.conf
usr/
usr/bin/
usr/bin/motion
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/motion.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/motion.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/motion.conf
usr/share/
usr/share/doc/
usr/share/doc/motion/
usr/share/doc/motion/CREDITS
usr/share/doc/motion/mask1.png
usr/share/doc/motion/motion_build.html
usr/share/doc/motion/motion_config.html
usr/share/doc/motion/motion_guide.html
usr/share/doc/motion/motion_stylesheet.css
usr/share/doc/motion/normal.jpg
usr/share/doc/motion/outputmotion1.jpg
usr/share/doc/motion/outputnormal1.jpg
usr/share/locale/
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/motion.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/motion.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/motion.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/motion.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/motion.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/motion.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/motion.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/motion.mo
usr/share/locale/li/
usr/share/locale/li/LC_MESSAGES/
usr/share/locale/li/LC_MESSAGES/motion.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/motion.mo
usr/share/locale/no/
usr/share/locale/no/LC_MESSAGES/
usr/share/locale/no/LC_MESSAGES/motion.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/motion.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/motion.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/motion.mo
usr/share/locale/zh/
usr/share/locale/zh/LC_MESSAGES/
usr/share/locale/zh/LC_MESSAGES/motion.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/motion.1.gz
