# extracted/gortr-0.15.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/gortr
usr/
usr/bin/
usr/bin/gortr
usr/bin/rtrdump
usr/bin/rtrmon
usr/lib/
usr/lib/gortr/
usr/lib/gortr/cf.pub
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gortr.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/gortr.conf
usr/share/
usr/share/licenses/
usr/share/licenses/gortr/
usr/share/licenses/gortr/LICENSE
