# extracted/libiio-0.26-2/files
%FILES%
usr/
usr/bin/
usr/bin/iio_attr
usr/bin/iio_genxml
usr/bin/iio_info
usr/bin/iio_readdev
usr/bin/iio_reg
usr/bin/iio_stresstest
usr/bin/iio_writedev
usr/bin/iiod
usr/include/
usr/include/iio.h
usr/lib/
usr/lib/libiio.so
usr/lib/libiio.so.0
usr/lib/libiio.so.0.26
usr/lib/pkgconfig/
usr/lib/pkgconfig/libiio.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/iiod.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/90-libiio.rules
usr/share/
usr/share/doc/
usr/share/doc/libiio/
usr/share/doc/libiio/CONTRIBUTING.md
usr/share/doc/libiio/Contributors.md
usr/share/doc/libiio/README.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/iio_attr.1.gz
usr/share/man/man1/iio_genxml.1.gz
usr/share/man/man1/iio_info.1.gz
usr/share/man/man1/iio_readdev.1.gz
usr/share/man/man1/iio_reg.1.gz
usr/share/man/man1/iio_stresstest.1.gz
usr/share/man/man1/iio_writedev.1.gz
usr/share/man/man1/iiod.1.gz
usr/share/man/man3/
usr/share/man/man3/libiio.3.gz
