# extracted/munin-node-2.0.76-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/munin-node
etc/munin/
etc/munin/munin-node.conf
etc/munin/plugin-conf.d/
etc/munin/plugins/
usr/
usr/bin/
usr/bin/munin-get
usr/bin/munin-node
usr/bin/munin-node-configure
usr/bin/munin-run
usr/bin/munindoc
usr/lib/
usr/lib/munin/
usr/lib/munin/munin-async
usr/lib/munin/munin-asyncd
usr/lib/munin/plugins/
usr/lib/munin/plugins/acpi
usr/lib/munin/plugins/amavis
usr/lib/munin/plugins/apache_accesses
usr/lib/munin/plugins/apache_processes
usr/lib/munin/plugins/apache_volume
usr/lib/munin/plugins/apc_envunit_
usr/lib/munin/plugins/apc_nis
usr/lib/munin/plugins/apt
usr/lib/munin/plugins/apt_all
usr/lib/munin/plugins/bind9
usr/lib/munin/plugins/bind9_rndc
usr/lib/munin/plugins/bonding_err_
usr/lib/munin/plugins/buddyinfo
usr/lib/munin/plugins/cmc_tc_sensor_
usr/lib/munin/plugins/courier_
usr/lib/munin/plugins/courier_mta_mailqueue
usr/lib/munin/plugins/courier_mta_mailstats
usr/lib/munin/plugins/courier_mta_mailvolume
usr/lib/munin/plugins/cps_
usr/lib/munin/plugins/cpu
usr/lib/munin/plugins/cpuspeed
usr/lib/munin/plugins/cupsys_pages
usr/lib/munin/plugins/df
usr/lib/munin/plugins/df_abs
usr/lib/munin/plugins/df_inode
usr/lib/munin/plugins/dhcpd3
usr/lib/munin/plugins/digitemp_
usr/lib/munin/plugins/diskstat_
usr/lib/munin/plugins/diskstats
usr/lib/munin/plugins/ejabberd_
usr/lib/munin/plugins/entropy
usr/lib/munin/plugins/exim_mailqueue
usr/lib/munin/plugins/exim_mailqueue_alt
usr/lib/munin/plugins/exim_mailstats
usr/lib/munin/plugins/fail2ban
usr/lib/munin/plugins/foldingathome
usr/lib/munin/plugins/foldingathome_rank
usr/lib/munin/plugins/foldingathome_wu
usr/lib/munin/plugins/forks
usr/lib/munin/plugins/freeradius_acct
usr/lib/munin/plugins/freeradius_auth
usr/lib/munin/plugins/freeradius_proxy_acct
usr/lib/munin/plugins/freeradius_proxy_auth
usr/lib/munin/plugins/fw_conntrack
usr/lib/munin/plugins/fw_forwarded_local
usr/lib/munin/plugins/fw_packets
usr/lib/munin/plugins/haproxy_
usr/lib/munin/plugins/haproxy_ng
usr/lib/munin/plugins/hddtemp
usr/lib/munin/plugins/hddtemp2
usr/lib/munin/plugins/hddtemp_smartctl
usr/lib/munin/plugins/hddtempd
usr/lib/munin/plugins/hp2000_
usr/lib/munin/plugins/http_loadtime
usr/lib/munin/plugins/if_
usr/lib/munin/plugins/if_err_
usr/lib/munin/plugins/ifx_concurrent_sessions_
usr/lib/munin/plugins/interrupts
usr/lib/munin/plugins/iostat
usr/lib/munin/plugins/iostat_ios
usr/lib/munin/plugins/ip_
usr/lib/munin/plugins/ipac-ng
usr/lib/munin/plugins/ipmi_
usr/lib/munin/plugins/ipmi_sensor_
usr/lib/munin/plugins/ircu
usr/lib/munin/plugins/irqstats
usr/lib/munin/plugins/load
usr/lib/munin/plugins/loggrep
usr/lib/munin/plugins/lpar_cpu
usr/lib/munin/plugins/lpstat
usr/lib/munin/plugins/mailman
usr/lib/munin/plugins/mailscanner
usr/lib/munin/plugins/mbmon_
usr/lib/munin/plugins/memcached_
usr/lib/munin/plugins/meminfo
usr/lib/munin/plugins/memory
usr/lib/munin/plugins/mhttping
usr/lib/munin/plugins/multiping
usr/lib/munin/plugins/multips
usr/lib/munin/plugins/multips_memory
usr/lib/munin/plugins/munin_stats
usr/lib/munin/plugins/munin_update
usr/lib/munin/plugins/mysql_
usr/lib/munin/plugins/mysql_bytes
usr/lib/munin/plugins/mysql_innodb
usr/lib/munin/plugins/mysql_isam_space_
usr/lib/munin/plugins/mysql_queries
usr/lib/munin/plugins/mysql_slowqueries
usr/lib/munin/plugins/mysql_threads
usr/lib/munin/plugins/named
usr/lib/munin/plugins/netopia
usr/lib/munin/plugins/netstat
usr/lib/munin/plugins/netstat_multi
usr/lib/munin/plugins/nfs4_client
usr/lib/munin/plugins/nfs_client
usr/lib/munin/plugins/nfsd
usr/lib/munin/plugins/nfsd4
usr/lib/munin/plugins/nginx_request
usr/lib/munin/plugins/nginx_status
usr/lib/munin/plugins/nomadix_users_
usr/lib/munin/plugins/ntp_
usr/lib/munin/plugins/ntp_kernel_err
usr/lib/munin/plugins/ntp_kernel_pll_freq
usr/lib/munin/plugins/ntp_kernel_pll_off
usr/lib/munin/plugins/ntp_offset
usr/lib/munin/plugins/ntp_states
usr/lib/munin/plugins/nut_misc
usr/lib/munin/plugins/nut_volts
usr/lib/munin/plugins/nutups_
usr/lib/munin/plugins/nvidia_
usr/lib/munin/plugins/open_files
usr/lib/munin/plugins/open_inodes
usr/lib/munin/plugins/openvpn
usr/lib/munin/plugins/perdition
usr/lib/munin/plugins/pgbouncer_connections
usr/lib/munin/plugins/pgbouncer_requests
usr/lib/munin/plugins/ping_
usr/lib/munin/plugins/plugin.sh
usr/lib/munin/plugins/plugins.history
usr/lib/munin/plugins/pm3users_
usr/lib/munin/plugins/pop_stats
usr/lib/munin/plugins/port_
usr/lib/munin/plugins/postfix_mailqueue
usr/lib/munin/plugins/postfix_mailstats
usr/lib/munin/plugins/postfix_mailvolume
usr/lib/munin/plugins/postgres_autovacuum
usr/lib/munin/plugins/postgres_bgwriter
usr/lib/munin/plugins/postgres_cache_
usr/lib/munin/plugins/postgres_checkpoints
usr/lib/munin/plugins/postgres_connections_
usr/lib/munin/plugins/postgres_connections_db
usr/lib/munin/plugins/postgres_locks_
usr/lib/munin/plugins/postgres_oldest_prepared_xact_
usr/lib/munin/plugins/postgres_prepared_xacts_
usr/lib/munin/plugins/postgres_querylength_
usr/lib/munin/plugins/postgres_scans_
usr/lib/munin/plugins/postgres_size_
usr/lib/munin/plugins/postgres_streaming_
usr/lib/munin/plugins/postgres_transactions_
usr/lib/munin/plugins/postgres_tuples_
usr/lib/munin/plugins/postgres_users
usr/lib/munin/plugins/postgres_xlog
usr/lib/munin/plugins/proc
usr/lib/munin/plugins/proc_pri
usr/lib/munin/plugins/processes
usr/lib/munin/plugins/proxy_plugin
usr/lib/munin/plugins/ps_
usr/lib/munin/plugins/psu_
usr/lib/munin/plugins/qmailqstat
usr/lib/munin/plugins/qmailscan
usr/lib/munin/plugins/qmailscan-simple
usr/lib/munin/plugins/quota_usage_
usr/lib/munin/plugins/samba
usr/lib/munin/plugins/selinux_avcstat
usr/lib/munin/plugins/sendmail_mailqueue
usr/lib/munin/plugins/sendmail_mailstats
usr/lib/munin/plugins/sendmail_mailtraffic
usr/lib/munin/plugins/sensors_
usr/lib/munin/plugins/slapd_
usr/lib/munin/plugins/slapd_bdb_cache_
usr/lib/munin/plugins/slony_lag_
usr/lib/munin/plugins/smart_
usr/lib/munin/plugins/snmp__cpuload
usr/lib/munin/plugins/snmp__df
usr/lib/munin/plugins/snmp__df_ram
usr/lib/munin/plugins/snmp__fc_if_
usr/lib/munin/plugins/snmp__fc_if_err_
usr/lib/munin/plugins/snmp__if_
usr/lib/munin/plugins/snmp__if_err_
usr/lib/munin/plugins/snmp__if_multi
usr/lib/munin/plugins/snmp__load
usr/lib/munin/plugins/snmp__memory
usr/lib/munin/plugins/snmp__netapp_diskusage_
usr/lib/munin/plugins/snmp__netapp_inodeusage_
usr/lib/munin/plugins/snmp__netstat
usr/lib/munin/plugins/snmp__print_pages
usr/lib/munin/plugins/snmp__print_supplies
usr/lib/munin/plugins/snmp__processes
usr/lib/munin/plugins/snmp__rdp_users
usr/lib/munin/plugins/snmp__sensors_fsc_bx_fan
usr/lib/munin/plugins/snmp__sensors_fsc_bx_temp
usr/lib/munin/plugins/snmp__sensors_fsc_fan
usr/lib/munin/plugins/snmp__sensors_fsc_temp
usr/lib/munin/plugins/snmp__sensors_mbm_fan
usr/lib/munin/plugins/snmp__sensors_mbm_temp
usr/lib/munin/plugins/snmp__sensors_mbm_volt
usr/lib/munin/plugins/snmp__swap
usr/lib/munin/plugins/snmp__uptime
usr/lib/munin/plugins/snmp__users
usr/lib/munin/plugins/snmp__winload
usr/lib/munin/plugins/snmp__winmem
usr/lib/munin/plugins/snort_alerts
usr/lib/munin/plugins/snort_bytes_pkt
usr/lib/munin/plugins/snort_drop_rate
usr/lib/munin/plugins/snort_pattern_match
usr/lib/munin/plugins/snort_pkts
usr/lib/munin/plugins/snort_traffic
usr/lib/munin/plugins/spamstats
usr/lib/munin/plugins/squeezebox_
usr/lib/munin/plugins/squid_cache
usr/lib/munin/plugins/squid_icp
usr/lib/munin/plugins/squid_objectsize
usr/lib/munin/plugins/squid_requests
usr/lib/munin/plugins/squid_traffic
usr/lib/munin/plugins/surfboard
usr/lib/munin/plugins/swap
usr/lib/munin/plugins/sybase_space
usr/lib/munin/plugins/tcp
usr/lib/munin/plugins/threads
usr/lib/munin/plugins/tomcat_
usr/lib/munin/plugins/tomcat_access
usr/lib/munin/plugins/tomcat_jvm
usr/lib/munin/plugins/tomcat_threads
usr/lib/munin/plugins/tomcat_volume
usr/lib/munin/plugins/uptime
usr/lib/munin/plugins/users
usr/lib/munin/plugins/varnish_
usr/lib/munin/plugins/vlan_
usr/lib/munin/plugins/vlan_inetuse_
usr/lib/munin/plugins/vlan_linkuse_
usr/lib/munin/plugins/vmstat
usr/lib/munin/plugins/vserver_cpu_
usr/lib/munin/plugins/vserver_loadavg
usr/lib/munin/plugins/vserver_resources
usr/lib/munin/plugins/yum
usr/lib/munin/plugins/zimbra_
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/munin-async.service
usr/lib/systemd/system/munin-node.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/munin-node.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/munin-node.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/munin-get.1p.gz
usr/share/man/man1/munin-node-configure.1p.gz
usr/share/man/man1/munin-node.1p.gz
usr/share/man/man1/munin-run.1p.gz
usr/share/man/man1/munindoc.1p.gz
usr/share/man/man3/
usr/share/man/man3/Munin::Common::Config.3pm.gz
usr/share/man/man3/Munin::Common::Daemon.3pm.gz
usr/share/man/man3/Munin::Common::Defaults.3pm.gz
usr/share/man/man3/Munin::Common::TLS.3pm.gz
usr/share/man/man3/Munin::Common::TLSClient.3pm.gz
usr/share/man/man3/Munin::Common::TLSServer.3pm.gz
usr/share/man/man3/Munin::Common::Timeout.3pm.gz
usr/share/man/man3/Munin::Node::Config.3pm.gz
usr/share/man/man3/Munin::Node::Configure::Debug.3pm.gz
usr/share/man/man3/Munin::Node::Configure::History.3pm.gz
usr/share/man/man3/Munin::Node::Configure::HostEnumeration.3pm.gz
usr/share/man/man3/Munin::Node::Configure::Plugin.3pm.gz
usr/share/man/man3/Munin::Node::Configure::PluginList.3pm.gz
usr/share/man/man3/Munin::Node::Logger.3pm.gz
usr/share/man/man3/Munin::Node::OS.3pm.gz
usr/share/man/man3/Munin::Node::SNMPConfig.3pm.gz
usr/share/man/man3/Munin::Node::Server.3pm.gz
usr/share/man/man3/Munin::Node::Service.3pm.gz
usr/share/man/man3/Munin::Node::Session.3pm.gz
usr/share/man/man3/Munin::Node::SpoolReader.3pm.gz
usr/share/man/man3/Munin::Node::SpoolWriter.3pm.gz
usr/share/man/man3/Munin::Node::Utils.3pm.gz
usr/share/man/man3/Munin::Plugin.3pm.gz
usr/share/man/man3/Munin::Plugin::Pgsql.3pm.gz
usr/share/man/man3/Munin::Plugin::SNMP.3pm.gz
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/Munin/
usr/share/perl5/vendor_perl/Munin/Common/
usr/share/perl5/vendor_perl/Munin/Common/Config.pm
usr/share/perl5/vendor_perl/Munin/Common/Daemon.pm
usr/share/perl5/vendor_perl/Munin/Common/Defaults.pm
usr/share/perl5/vendor_perl/Munin/Common/DictFile.pm
usr/share/perl5/vendor_perl/Munin/Common/SyncDictFile.pm
usr/share/perl5/vendor_perl/Munin/Common/TLS.pm
usr/share/perl5/vendor_perl/Munin/Common/TLSClient.pm
usr/share/perl5/vendor_perl/Munin/Common/TLSServer.pm
usr/share/perl5/vendor_perl/Munin/Common/Timeout.pm
usr/share/perl5/vendor_perl/Munin/Node/
usr/share/perl5/vendor_perl/Munin/Node/Config.pm
usr/share/perl5/vendor_perl/Munin/Node/Configure/
usr/share/perl5/vendor_perl/Munin/Node/Configure/Debug.pm
usr/share/perl5/vendor_perl/Munin/Node/Configure/History.pm
usr/share/perl5/vendor_perl/Munin/Node/Configure/HostEnumeration.pm
usr/share/perl5/vendor_perl/Munin/Node/Configure/Plugin.pm
usr/share/perl5/vendor_perl/Munin/Node/Configure/PluginList.pm
usr/share/perl5/vendor_perl/Munin/Node/Logger.pm
usr/share/perl5/vendor_perl/Munin/Node/OS.pm
usr/share/perl5/vendor_perl/Munin/Node/SNMPConfig.pm
usr/share/perl5/vendor_perl/Munin/Node/Server.pm
usr/share/perl5/vendor_perl/Munin/Node/Service.pm
usr/share/perl5/vendor_perl/Munin/Node/Session.pm
usr/share/perl5/vendor_perl/Munin/Node/SpoolReader.pm
usr/share/perl5/vendor_perl/Munin/Node/SpoolWriter.pm
usr/share/perl5/vendor_perl/Munin/Node/Utils.pm
usr/share/perl5/vendor_perl/Munin/Plugin.pm
usr/share/perl5/vendor_perl/Munin/Plugin/
usr/share/perl5/vendor_perl/Munin/Plugin/Pgsql.pm
usr/share/perl5/vendor_perl/Munin/Plugin/SNMP.pm
var/
var/lib/
var/lib/munin/
var/lib/munin/plugin-state/
var/lib/munin/spool/
var/log/
var/log/munin/
