# extracted/shairplay-20180824.096b61a-5/files
%FILES%
etc/
etc/shairplay/
etc/shairplay/airport.key
usr/
usr/bin/
usr/bin/shairplay
usr/include/
usr/include/shairplay/
usr/include/shairplay/dnssd.h
usr/include/shairplay/raop.h
usr/lib/
usr/lib/libshairplay.so
usr/lib/libshairplay.so.0
usr/lib/libshairplay.so.0.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shairplay.service
usr/share/
usr/share/doc/
usr/share/doc/shairplay/
usr/share/doc/shairplay/README.md
usr/share/licenses/
usr/share/licenses/shairplay/
usr/share/licenses/shairplay/LICENSE
