# extracted/munin-2.0.76-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/munin
etc/munin/
etc/munin/munin-conf.d/
etc/munin/munin-cron-entry
etc/munin/munin.conf
etc/munin/static/
etc/munin/static/definitions.html
etc/munin/static/dynazoom.html
etc/munin/static/favicon.ico
etc/munin/static/formatdate.js
etc/munin/static/logo-h.png
etc/munin/static/logo.png
etc/munin/static/querystring.js
etc/munin/static/style-1.2.css
etc/munin/static/style-new.css
etc/munin/static/style.css
etc/munin/static/zoom.js
etc/munin/templates/
etc/munin/templates/munin-categoryview.tmpl
etc/munin/templates/munin-comparison-day.tmpl
etc/munin/templates/munin-comparison-month.tmpl
etc/munin/templates/munin-comparison-week.tmpl
etc/munin/templates/munin-comparison-year.tmpl
etc/munin/templates/munin-domainview.tmpl
etc/munin/templates/munin-dynazoom.tmpl
etc/munin/templates/munin-nodeview.tmpl
etc/munin/templates/munin-overview.tmpl
etc/munin/templates/munin-problemview.tmpl
etc/munin/templates/munin-serviceview.tmpl
etc/munin/templates/partial/
etc/munin/templates/partial/bottom_navigation.tmpl
etc/munin/templates/partial/footer.tmpl
etc/munin/templates/partial/generated_by.tmpl
etc/munin/templates/partial/head.tmpl
etc/munin/templates/partial/logo_navigation.tmpl
etc/munin/templates/partial/logo_navigation_comparison.tmpl
etc/munin/templates/partial/logo_navigation_problem.tmpl
etc/munin/templates/partial/logo_path.tmpl
etc/munin/templates/partial/navigation.tmpl
etc/munin/templates/partial/path.tmpl
usr/
usr/bin/
usr/bin/munin-check
usr/bin/munin-cron
usr/lib/
usr/lib/munin/
usr/lib/munin/DejaVuSans.ttf
usr/lib/munin/DejaVuSansMono.ttf
usr/lib/munin/munin-datafile2storable
usr/lib/munin/munin-graph
usr/lib/munin/munin-html
usr/lib/munin/munin-limits
usr/lib/munin/munin-storable2datafile
usr/lib/munin/munin-update
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/munin-graph.service
usr/lib/systemd/system/munin-graph.socket
usr/lib/systemd/system/munin-html.service
usr/lib/systemd/system/munin-html.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/munin.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/munin.conf
usr/share/
usr/share/fontconfig/
usr/share/fontconfig/conf.avail/
usr/share/fontconfig/conf.avail/08-munin-font-dir.conf
usr/share/fontconfig/conf.default/
usr/share/fontconfig/conf.default/08-munin-font-dir.conf
usr/share/man/
usr/share/man/man3/
usr/share/man/man3/Munin::Master::Config.3pm.gz
usr/share/man/man3/Munin::Master::Group.3pm.gz
usr/share/man/man3/Munin::Master::GroupRepository.3pm.gz
usr/share/man/man3/Munin::Master::HTMLOld.3pm.gz
usr/share/man/man3/Munin::Master::Host.3pm.gz
usr/share/man/man3/Munin::Master::LimitsOld.3pm.gz
usr/share/man/man3/Munin::Master::Logger.3pm.gz
usr/share/man/man3/Munin::Master::Node.3pm.gz
usr/share/man/man3/Munin::Master::ProcessManager.3pm.gz
usr/share/man/man3/Munin::Master::Update.3pm.gz
usr/share/man/man3/Munin::Master::UpdateWorker.3pm.gz
usr/share/man/man3/Munin::Master::Utils.3pm.gz
usr/share/man/man3/Munin::Master::Worker.3pm.gz
usr/share/munin/
usr/share/munin/cgi/
usr/share/munin/cgi/munin-cgi-graph
usr/share/munin/cgi/munin-cgi-html
usr/share/munin/www/
usr/share/munin/www/.htaccess
usr/share/perl5/
usr/share/perl5/vendor_perl/
usr/share/perl5/vendor_perl/Munin/
usr/share/perl5/vendor_perl/Munin/Master/
usr/share/perl5/vendor_perl/Munin/Master/Config.pm
usr/share/perl5/vendor_perl/Munin/Master/GraphOld.pm
usr/share/perl5/vendor_perl/Munin/Master/Group.pm
usr/share/perl5/vendor_perl/Munin/Master/GroupRepository.pm
usr/share/perl5/vendor_perl/Munin/Master/HTMLConfig.pm
usr/share/perl5/vendor_perl/Munin/Master/HTMLOld.pm
usr/share/perl5/vendor_perl/Munin/Master/Host.pm
usr/share/perl5/vendor_perl/Munin/Master/LimitsOld.pm
usr/share/perl5/vendor_perl/Munin/Master/Logger.pm
usr/share/perl5/vendor_perl/Munin/Master/Node.pm
usr/share/perl5/vendor_perl/Munin/Master/ProcessManager.pm
usr/share/perl5/vendor_perl/Munin/Master/Update.pm
usr/share/perl5/vendor_perl/Munin/Master/UpdateWorker.pm
usr/share/perl5/vendor_perl/Munin/Master/Utils.pm
usr/share/perl5/vendor_perl/Munin/Master/Worker.pm
var/
var/lib/
var/lib/munin/
var/lib/munin/cgi-tmp/
var/lib/munin/spool/
var/log/
var/log/munin/
