# extracted/open-isns-0.103-1/files
%FILES%
etc/
etc/isns/
etc/isns/isnsadm.conf
etc/isns/isnsd.conf
etc/isns/isnsdd.conf
usr/
usr/bin/
usr/bin/isnsadm
usr/bin/isnsd
usr/bin/isnsdd
usr/include/
usr/include/libisns/
usr/include/libisns/attrs.h
usr/include/libisns/buffer.h
usr/include/libisns/isns-proto.h
usr/include/libisns/isns.h
usr/include/libisns/message.h
usr/include/libisns/paths.h
usr/include/libisns/source.h
usr/include/libisns/types.h
usr/include/libisns/util.h
usr/lib/
usr/lib/libisns.so
usr/lib/libisns.so.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/libisns.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/isnsd.service
usr/lib/systemd/system/isnsd.socket
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/isns_config.5.gz
usr/share/man/man8/
usr/share/man/man8/isnsadm.8.gz
usr/share/man/man8/isnsd.8.gz
usr/share/man/man8/isnsdd.8.gz
var/
var/lib/
var/lib/isns/
