# extracted/atftp-0.8.0-4/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/atftpd
usr/
usr/bin/
usr/bin/atftp
usr/bin/atftpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/atftpd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/atftp.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/atftp.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/atftp.1.gz
usr/share/man/man8/
usr/share/man/man8/atftpd.8.gz
