# extracted/deepin-application-manager-1.2.5-1/files
%FILES%
usr/
usr/bin/
usr/bin/app-identifier
usr/bin/dde-am
usr/bin/dde-application-manager
usr/lib/
usr/lib/deepin/
usr/lib/deepin/application-manager/
usr/lib/deepin/application-manager/app-launch-helper
usr/lib/deepin/application-manager/app-update-notifier
usr/lib/deepin/application-manager/dockEnv.sh
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/org.desktopspec.ApplicationUpdateNotifier1.service
usr/lib/systemd/user/
usr/lib/systemd/user/dde-session-initialized.target.wants/
usr/lib/systemd/user/dde-session-initialized.target.wants/org.desktopspec.ApplicationManager1.service
usr/lib/systemd/user/org.desktopspec.ApplicationManager1.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.desktopspec.ApplicationManager1.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.desktopspec.ApplicationUpdateNotifier1.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.desktopspec.ApplicationUpdateNotifier1.conf
usr/share/dde-application-manager/
usr/share/dde-application-manager/org.desktopspec.ApplicationManager1.Application.xml
usr/share/dde-application-manager/org.desktopspec.ApplicationManager1.Instance.xml
usr/share/dde-application-manager/org.desktopspec.ApplicationManager1.xml
usr/share/dde-application-manager/org.desktopspec.JobManager1.Job.xml
usr/share/dde-application-manager/org.desktopspec.JobManager1.xml
usr/share/dde-application-manager/org.desktopspec.MimeManager1.xml
usr/share/dde-application-manager/org.desktopspec.ObjectManager1.xml
usr/share/deepin/
usr/share/deepin/dde-application-manager/
usr/share/deepin/dde-application-manager/hooks.d/
usr/share/deepin/dde-application-manager/hooks.d/1-dockEnv.json
usr/share/dsg/
usr/share/dsg/configs/
usr/share/dsg/configs/dde-application-manager/
usr/share/dsg/configs/dde-application-manager/com.deepin.dde.appearance.json
usr/share/dsg/configs/dde-application-manager/com.deepin.dde.dock.json
usr/share/dsg/configs/dde-application-manager/com.deepin.dde.launcher.json
usr/share/dsg/configs/dde-application-manager/com.deepin.dde.startdde.json
usr/share/dsg/configs/dde-application-manager/com.deepin.xsettings.json
usr/share/dsg/configs/org.deepin.dde.application-manager/
usr/share/dsg/configs/org.deepin.dde.application-manager/org.deepin.dde.am.json
