# extracted/bftpd-6.1-4/files
%FILES%
etc/
etc/bftpd.conf
etc/logrotate.d/
etc/logrotate.d/bftpd
etc/pam.d/
etc/pam.d/bftpd
usr/
usr/bin/
usr/bin/bftpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bftpd.service
usr/lib/systemd/system/bftpd.socket
usr/lib/systemd/system/bftpd@.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/bftpd.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/bftpd.8.gz
