# extracted/irker-2.24-1/files
%FILES%
etc/
etc/irker/
etc/irker/irkerd.env
usr/
usr/bin/
usr/bin/irk
usr/bin/irkerd
usr/bin/irkerhook.py
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/irkerd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/irker.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/irker.conf
usr/share/
usr/share/doc/
usr/share/doc/irker/
usr/share/doc/irker/NEWS
usr/share/doc/irker/README
usr/share/doc/irker/hacking.adoc
usr/share/doc/irker/security.adoc
usr/share/licenses/
usr/share/licenses/irker/
usr/share/licenses/irker/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/irk.1.gz
usr/share/man/man1/irkerhook.1.gz
usr/share/man/man8/
usr/share/man/man8/irkerd.8.gz
