# extracted/kube-apiserver-1.31.3-1/files
%FILES%
etc/
etc/kubernetes/
etc/kubernetes/kube-apiserver.env
usr/
usr/bin/
usr/bin/kube-apiserver
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/kube-apiserver.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/kube-apiserver.1.gz
