# extracted/mariadb-11.6.2-2/files
%FILES%
etc/
etc/my.cnf
etc/my.cnf.d/
etc/my.cnf.d/client.cnf
etc/my.cnf.d/enable_encryption.preset
etc/my.cnf.d/hashicorp_key_management.cnf
etc/my.cnf.d/mysql-clients.cnf
etc/my.cnf.d/provider_bzip2.cnf
etc/my.cnf.d/provider_lz4.cnf
etc/my.cnf.d/provider_lzma.cnf
etc/my.cnf.d/s3.cnf
etc/my.cnf.d/server.cnf
etc/my.cnf.d/spider.cnf
etc/security/
etc/security/user_map.conf
usr/
usr/bin/
usr/bin/aria_chk
usr/bin/aria_dump_log
usr/bin/aria_ftdump
usr/bin/aria_pack
usr/bin/aria_read_log
usr/bin/aria_s3_copy
usr/bin/galera_new_cluster
usr/bin/galera_recovery
usr/bin/innochecksum
usr/bin/mariabackup
usr/bin/mariadb-access
usr/bin/mariadb-backup
usr/bin/mariadb-client-test
usr/bin/mariadb-client-test-embedded
usr/bin/mariadb-convert-table-format
usr/bin/mariadb-dumpslow
usr/bin/mariadb-embedded
usr/bin/mariadb-find-rows
usr/bin/mariadb-fix-extensions
usr/bin/mariadb-hotcopy
usr/bin/mariadb-install-db
usr/bin/mariadb-ldb
usr/bin/mariadb-secure-installation
usr/bin/mariadb-service-convert
usr/bin/mariadb-setpermission
usr/bin/mariadb-test-embedded
usr/bin/mariadb-tzinfo-to-sql
usr/bin/mariadb-waitpid
usr/bin/mariadbd
usr/bin/mariadbd-multi
usr/bin/mariadbd-safe
usr/bin/mariadbd-safe-helper
usr/bin/mbstream
usr/bin/msql2mysql
usr/bin/my_print_defaults
usr/bin/myisam_ftdump
usr/bin/myisamchk
usr/bin/myisamlog
usr/bin/myisampack
usr/bin/myrocks_hotbackup
usr/bin/mysql_client_test
usr/bin/mysql_client_test_embedded
usr/bin/mysql_convert_table_format
usr/bin/mysql_embedded
usr/bin/mysql_find_rows
usr/bin/mysql_fix_extensions
usr/bin/mysql_install_db
usr/bin/mysql_ldb
usr/bin/mysql_secure_installation
usr/bin/mysql_setpermission
usr/bin/mysql_tzinfo_to_sql
usr/bin/mysql_waitpid
usr/bin/mysqlaccess
usr/bin/mysqld
usr/bin/mysqld_multi
usr/bin/mysqld_safe
usr/bin/mysqld_safe_helper
usr/bin/mysqldumpslow
usr/bin/mysqlhotcopy
usr/bin/mysqltest_embedded
usr/bin/perror
usr/bin/replace
usr/bin/resolve_stack_dump
usr/bin/resolveip
usr/bin/sst_dump
usr/bin/test-connect-t
usr/bin/wsrep_sst_backup
usr/bin/wsrep_sst_common
usr/bin/wsrep_sst_mariabackup
usr/bin/wsrep_sst_mysqldump
usr/bin/wsrep_sst_rsync
usr/bin/wsrep_sst_rsync_wan
usr/lib/
usr/lib/mysql/
usr/lib/mysql/plugin/
usr/lib/mysql/plugin/adt_null.so
usr/lib/mysql/plugin/auth_0x0100.so
usr/lib/mysql/plugin/auth_ed25519.so
usr/lib/mysql/plugin/auth_gssapi.so
usr/lib/mysql/plugin/auth_pam.so
usr/lib/mysql/plugin/auth_pam_tool_dir/
usr/lib/mysql/plugin/auth_pam_tool_dir/auth_pam_tool
usr/lib/mysql/plugin/auth_pam_v1.so
usr/lib/mysql/plugin/auth_parsec.so
usr/lib/mysql/plugin/auth_test_plugin.so
usr/lib/mysql/plugin/cracklib_password_check.so
usr/lib/mysql/plugin/daemon_example.ini
usr/lib/mysql/plugin/debug_key_management.so
usr/lib/mysql/plugin/dialog_examples.so
usr/lib/mysql/plugin/disks.so
usr/lib/mysql/plugin/example_key_management.so
usr/lib/mysql/plugin/file_key_management.so
usr/lib/mysql/plugin/func_test.so
usr/lib/mysql/plugin/ha_archive.so
usr/lib/mysql/plugin/ha_blackhole.so
usr/lib/mysql/plugin/ha_connect.so
usr/lib/mysql/plugin/ha_federatedx.so
usr/lib/mysql/plugin/ha_mroonga.so
usr/lib/mysql/plugin/ha_oqgraph.so
usr/lib/mysql/plugin/ha_rocksdb.so
usr/lib/mysql/plugin/ha_s3.so
usr/lib/mysql/plugin/ha_sphinx.so
usr/lib/mysql/plugin/ha_spider.so
usr/lib/mysql/plugin/ha_test_sql_discovery.so
usr/lib/mysql/plugin/handlersocket.so
usr/lib/mysql/plugin/hashicorp_key_management.so
usr/lib/mysql/plugin/libdaemon_example.so
usr/lib/mysql/plugin/locales.so
usr/lib/mysql/plugin/metadata_lock_info.so
usr/lib/mysql/plugin/mypluglib.so
usr/lib/mysql/plugin/password_reuse_check.so
usr/lib/mysql/plugin/provider_bzip2.so
usr/lib/mysql/plugin/provider_lz4.so
usr/lib/mysql/plugin/provider_lzma.so
usr/lib/mysql/plugin/qa_auth_client.so
usr/lib/mysql/plugin/qa_auth_interface.so
usr/lib/mysql/plugin/qa_auth_server.so
usr/lib/mysql/plugin/query_cache_info.so
usr/lib/mysql/plugin/query_response_time.so
usr/lib/mysql/plugin/server_audit.so
usr/lib/mysql/plugin/simple_password_check.so
usr/lib/mysql/plugin/sql_errlog.so
usr/lib/mysql/plugin/test_sql_service.so
usr/lib/mysql/plugin/test_versioning.so
usr/lib/mysql/plugin/type_mysql_json.so
usr/lib/mysql/plugin/type_mysql_timestamp.so
usr/lib/mysql/plugin/type_test.so
usr/lib/mysql/plugin/wsrep_info.so
usr/lib/security/
usr/lib/security/pam_user_map.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mariadb-extra.socket
usr/lib/systemd/system/mariadb-extra@.socket
usr/lib/systemd/system/mariadb.service
usr/lib/systemd/system/mariadb.socket
usr/lib/systemd/system/mariadb@.service
usr/lib/systemd/system/mariadb@.socket
usr/lib/systemd/system/mysql.service
usr/lib/systemd/system/mysqld.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mariadb.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mariadb.conf
usr/share/
usr/share/doc/
usr/share/doc/mariadb/
usr/share/doc/mariadb/CREDITS
usr/share/doc/mariadb/INSTALL-BINARY
usr/share/doc/mariadb/README-wsrep
usr/share/doc/mariadb/README.md
usr/share/doc/mariadb/THIRDPARTY
usr/share/doc/mariadb/groonga-normalizer-mysql/
usr/share/doc/mariadb/groonga-normalizer-mysql/README.md
usr/share/doc/mariadb/groonga-normalizer-mysql/lgpl-2.0.txt
usr/share/doc/mariadb/groonga/
usr/share/doc/mariadb/groonga/COPYING
usr/share/doc/mariadb/groonga/README.md
usr/share/doc/mariadb/hashicorp_key_management.txt
usr/share/licenses/
usr/share/licenses/mariadb/
usr/share/licenses/mariadb/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/aria_chk.1.gz
usr/share/man/man1/aria_dump_log.1.gz
usr/share/man/man1/aria_ftdump.1.gz
usr/share/man/man1/aria_pack.1.gz
usr/share/man/man1/aria_read_log.1.gz
usr/share/man/man1/aria_s3_copy.1.gz
usr/share/man/man1/galera_new_cluster.1.gz
usr/share/man/man1/galera_recovery.1.gz
usr/share/man/man1/innochecksum.1.gz
usr/share/man/man1/mariabackup.1.gz
usr/share/man/man1/mariadb-access.1.gz
usr/share/man/man1/mariadb-backup.1.gz
usr/share/man/man1/mariadb-client-test-embedded.1.gz
usr/share/man/man1/mariadb-client-test.1.gz
usr/share/man/man1/mariadb-convert-table-format.1.gz
usr/share/man/man1/mariadb-dumpslow.1.gz
usr/share/man/man1/mariadb-embedded.1.gz
usr/share/man/man1/mariadb-find-rows.1.gz
usr/share/man/man1/mariadb-fix-extensions.1.gz
usr/share/man/man1/mariadb-hotcopy.1.gz
usr/share/man/man1/mariadb-install-db.1.gz
usr/share/man/man1/mariadb-ldb.1.gz
usr/share/man/man1/mariadb-secure-installation.1.gz
usr/share/man/man1/mariadb-service-convert.1.gz
usr/share/man/man1/mariadb-setpermission.1.gz
usr/share/man/man1/mariadb-test-embedded.1.gz
usr/share/man/man1/mariadb-tzinfo-to-sql.1.gz
usr/share/man/man1/mariadb-waitpid.1.gz
usr/share/man/man1/mariadbd-multi.1.gz
usr/share/man/man1/mariadbd-safe-helper.1.gz
usr/share/man/man1/mariadbd-safe.1.gz
usr/share/man/man1/mbstream.1.gz
usr/share/man/man1/msql2mysql.1.gz
usr/share/man/man1/my_print_defaults.1.gz
usr/share/man/man1/my_safe_process.1.gz
usr/share/man/man1/myisam_ftdump.1.gz
usr/share/man/man1/myisamchk.1.gz
usr/share/man/man1/myisamlog.1.gz
usr/share/man/man1/myisampack.1.gz
usr/share/man/man1/myrocks_hotbackup.1.gz
usr/share/man/man1/mysql-stress-test.pl.1.gz
usr/share/man/man1/mysql.server.1.gz
usr/share/man/man1/mysql_client_test.1.gz
usr/share/man/man1/mysql_client_test_embedded.1.gz
usr/share/man/man1/mysql_convert_table_format.1.gz
usr/share/man/man1/mysql_embedded.1.gz
usr/share/man/man1/mysql_find_rows.1.gz
usr/share/man/man1/mysql_fix_extensions.1.gz
usr/share/man/man1/mysql_install_db.1.gz
usr/share/man/man1/mysql_ldb.1.gz
usr/share/man/man1/mysql_secure_installation.1.gz
usr/share/man/man1/mysql_setpermission.1.gz
usr/share/man/man1/mysql_tzinfo_to_sql.1.gz
usr/share/man/man1/mysql_waitpid.1.gz
usr/share/man/man1/mysqlaccess.1.gz
usr/share/man/man1/mysqld_multi.1.gz
usr/share/man/man1/mysqld_safe.1.gz
usr/share/man/man1/mysqld_safe_helper.1.gz
usr/share/man/man1/mysqldumpslow.1.gz
usr/share/man/man1/mysqlhotcopy.1.gz
usr/share/man/man1/mysqltest_embedded.1.gz
usr/share/man/man1/mytop.1.gz
usr/share/man/man1/perror.1.gz
usr/share/man/man1/replace.1.gz
usr/share/man/man1/resolve_stack_dump.1.gz
usr/share/man/man1/resolveip.1.gz
usr/share/man/man1/wsrep_sst_backup.1.gz
usr/share/man/man1/wsrep_sst_common.1.gz
usr/share/man/man1/wsrep_sst_mariabackup.1.gz
usr/share/man/man1/wsrep_sst_mysqldump.1.gz
usr/share/man/man1/wsrep_sst_rsync.1.gz
usr/share/man/man1/wsrep_sst_rsync_wan.1.gz
usr/share/man/man3/
usr/share/man/man3/mariadb_cancel.3.gz
usr/share/man/man3/mariadb_connection.3.gz
usr/share/man/man3/mariadb_dyncol_check.3.gz
usr/share/man/man3/mariadb_dyncol_column_cmp_named.3.gz
usr/share/man/man3/mariadb_dyncol_column_count.3.gz
usr/share/man/man3/mariadb_dyncol_create_many_named.3.gz
usr/share/man/man3/mariadb_dyncol_create_many_num.3.gz
usr/share/man/man3/mariadb_dyncol_exists_named.3.gz
usr/share/man/man3/mariadb_dyncol_exists_num.3.gz
usr/share/man/man3/mariadb_dyncol_free.3.gz
usr/share/man/man3/mariadb_dyncol_list_named.3.gz
usr/share/man/man3/mariadb_dyncol_list_num.3.gz
usr/share/man/man3/mariadb_dyncol_unpack.3.gz
usr/share/man/man3/mariadb_dyncol_update_many_named.3.gz
usr/share/man/man3/mariadb_dyncol_update_many_num.3.gz
usr/share/man/man3/mariadb_get_infov.3.gz
usr/share/man/man3/mariadb_reconnect.3.gz
usr/share/man/man3/mariadb_rpl_close.3.gz
usr/share/man/man3/mariadb_rpl_fetch.3.gz
usr/share/man/man3/mariadb_rpl_get_optionsv.3.gz
usr/share/man/man3/mariadb_rpl_open.3.gz
usr/share/man/man3/mariadb_rpl_optionsv.3.gz
usr/share/man/man3/mariadb_stmt_execute_direct.3.gz
usr/share/man/man3/mariadb_stmt_fetch_fields.3.gz
usr/share/man/man3/mysql_affected_rows.3.gz
usr/share/man/man3/mysql_autocommit.3.gz
usr/share/man/man3/mysql_change_user.3.gz
usr/share/man/man3/mysql_close.3.gz
usr/share/man/man3/mysql_commit.3.gz
usr/share/man/man3/mysql_data_seek.3.gz
usr/share/man/man3/mysql_errno.3.gz
usr/share/man/man3/mysql_error.3.gz
usr/share/man/man3/mysql_fetch_field.3.gz
usr/share/man/man3/mysql_fetch_field_direct.3.gz
usr/share/man/man3/mysql_fetch_fields.3.gz
usr/share/man/man3/mysql_fetch_lengths.3.gz
usr/share/man/man3/mysql_fetch_row.3.gz
usr/share/man/man3/mysql_field_count.3.gz
usr/share/man/man3/mysql_field_seek.3.gz
usr/share/man/man3/mysql_field_tell.3.gz
usr/share/man/man3/mysql_free_result.3.gz
usr/share/man/man3/mysql_get_character_set_info.3.gz
usr/share/man/man3/mysql_get_client_info.3.gz
usr/share/man/man3/mysql_get_client_version.3.gz
usr/share/man/man3/mysql_get_host_info.3.gz
usr/share/man/man3/mysql_get_proto_info.3.gz
usr/share/man/man3/mysql_get_server_info.3.gz
usr/share/man/man3/mysql_get_server_version.3.gz
usr/share/man/man3/mysql_get_socket.3.gz
usr/share/man/man3/mysql_get_ssl_cipher.3.gz
usr/share/man/man3/mysql_hex_string.3.gz
usr/share/man/man3/mysql_info.3.gz
usr/share/man/man3/mysql_init.3.gz
usr/share/man/man3/mysql_kill.3.gz
usr/share/man/man3/mysql_more_results.3.gz
usr/share/man/man3/mysql_next_result.3.gz
usr/share/man/man3/mysql_num_fields.3.gz
usr/share/man/man3/mysql_num_rows.3.gz
usr/share/man/man3/mysql_options.3.gz
usr/share/man/man3/mysql_options4.3.gz
usr/share/man/man3/mysql_optionsv.3.gz
usr/share/man/man3/mysql_ping.3.gz
usr/share/man/man3/mysql_query.3.gz
usr/share/man/man3/mysql_read_query_result.3.gz
usr/share/man/man3/mysql_real_connect.3.gz
usr/share/man/man3/mysql_real_escape_string.3.gz
usr/share/man/man3/mysql_real_query.3.gz
usr/share/man/man3/mysql_refresh.3.gz
usr/share/man/man3/mysql_reset_connection.3.gz
usr/share/man/man3/mysql_rollback.3.gz
usr/share/man/man3/mysql_row_seek.3.gz
usr/share/man/man3/mysql_row_tell.3.gz
usr/share/man/man3/mysql_select_db.3.gz
usr/share/man/man3/mysql_send_query.3.gz
usr/share/man/man3/mysql_server_end.3.gz
usr/share/man/man3/mysql_server_init.3.gz
usr/share/man/man3/mysql_session_track_get_first.3.gz
usr/share/man/man3/mysql_session_track_get_next.3.gz
usr/share/man/man3/mysql_set_character_set.3.gz
usr/share/man/man3/mysql_set_server_option.3.gz
usr/share/man/man3/mysql_shutdown.3.gz
usr/share/man/man3/mysql_sqlstate.3.gz
usr/share/man/man3/mysql_ssl_set.3.gz
usr/share/man/man3/mysql_stat.3.gz
usr/share/man/man3/mysql_stmt_affected_rows.3.gz
usr/share/man/man3/mysql_stmt_attr_get.3.gz
usr/share/man/man3/mysql_stmt_attr_set.3.gz
usr/share/man/man3/mysql_stmt_bind_param.3.gz
usr/share/man/man3/mysql_stmt_bind_result.3.gz
usr/share/man/man3/mysql_stmt_close.3.gz
usr/share/man/man3/mysql_stmt_data_seek.3.gz
usr/share/man/man3/mysql_stmt_errno.3.gz
usr/share/man/man3/mysql_stmt_error.3.gz
usr/share/man/man3/mysql_stmt_execute.3.gz
usr/share/man/man3/mysql_stmt_fetch.3.gz
usr/share/man/man3/mysql_stmt_fetch_column.3.gz
usr/share/man/man3/mysql_stmt_field_count.3.gz
usr/share/man/man3/mysql_stmt_free_result.3.gz
usr/share/man/man3/mysql_stmt_init.3.gz
usr/share/man/man3/mysql_stmt_insert_id.3.gz
usr/share/man/man3/mysql_stmt_more_results.3.gz
usr/share/man/man3/mysql_stmt_next_result.3.gz
usr/share/man/man3/mysql_stmt_num_rows.3.gz
usr/share/man/man3/mysql_stmt_param_count.3.gz
usr/share/man/man3/mysql_stmt_param_metadata.3.gz
usr/share/man/man3/mysql_stmt_prepare.3.gz
usr/share/man/man3/mysql_stmt_reset.3.gz
usr/share/man/man3/mysql_stmt_result_metadata.3.gz
usr/share/man/man3/mysql_stmt_row_seek.3.gz
usr/share/man/man3/mysql_stmt_row_tell.3.gz
usr/share/man/man3/mysql_stmt_send_long_data.3.gz
usr/share/man/man3/mysql_stmt_sqlstate.3.gz
usr/share/man/man3/mysql_stmt_store_result.3.gz
usr/share/man/man3/mysql_stmt_warning_count.3.gz
usr/share/man/man3/mysql_store_result.3.gz
usr/share/man/man3/mysql_thread_end.3.gz
usr/share/man/man3/mysql_thread_id.3.gz
usr/share/man/man3/mysql_thread_init.3.gz
usr/share/man/man3/mysql_use_result.3.gz
usr/share/man/man3/mysql_warning_count.3.gz
usr/share/man/man8/
usr/share/man/man8/mariadbd.8.gz
usr/share/man/man8/mysqld.8.gz
usr/share/mysql/
usr/share/mysql/bulgarian/
usr/share/mysql/bulgarian/errmsg.sys
usr/share/mysql/charsets/
usr/share/mysql/charsets/Index.xml
usr/share/mysql/charsets/README
usr/share/mysql/charsets/armscii8.xml
usr/share/mysql/charsets/ascii.xml
usr/share/mysql/charsets/cp1250.xml
usr/share/mysql/charsets/cp1251.xml
usr/share/mysql/charsets/cp1256.xml
usr/share/mysql/charsets/cp1257.xml
usr/share/mysql/charsets/cp850.xml
usr/share/mysql/charsets/cp852.xml
usr/share/mysql/charsets/cp866.xml
usr/share/mysql/charsets/dec8.xml
usr/share/mysql/charsets/geostd8.xml
usr/share/mysql/charsets/greek.xml
usr/share/mysql/charsets/hebrew.xml
usr/share/mysql/charsets/hp8.xml
usr/share/mysql/charsets/keybcs2.xml
usr/share/mysql/charsets/koi8r.xml
usr/share/mysql/charsets/koi8u.xml
usr/share/mysql/charsets/latin1.xml
usr/share/mysql/charsets/latin2.xml
usr/share/mysql/charsets/latin5.xml
usr/share/mysql/charsets/latin7.xml
usr/share/mysql/charsets/macce.xml
usr/share/mysql/charsets/macroman.xml
usr/share/mysql/charsets/swe7.xml
usr/share/mysql/chinese/
usr/share/mysql/chinese/errmsg.sys
usr/share/mysql/czech/
usr/share/mysql/czech/errmsg.sys
usr/share/mysql/danish/
usr/share/mysql/danish/errmsg.sys
usr/share/mysql/dutch/
usr/share/mysql/dutch/errmsg.sys
usr/share/mysql/english/
usr/share/mysql/english/errmsg.sys
usr/share/mysql/estonian/
usr/share/mysql/estonian/errmsg.sys
usr/share/mysql/fill_help_tables.sql
usr/share/mysql/french/
usr/share/mysql/french/errmsg.sys
usr/share/mysql/georgian/
usr/share/mysql/georgian/errmsg.sys
usr/share/mysql/german/
usr/share/mysql/german/errmsg.sys
usr/share/mysql/greek/
usr/share/mysql/greek/errmsg.sys
usr/share/mysql/hindi/
usr/share/mysql/hindi/errmsg.sys
usr/share/mysql/hungarian/
usr/share/mysql/hungarian/errmsg.sys
usr/share/mysql/italian/
usr/share/mysql/italian/errmsg.sys
usr/share/mysql/japanese/
usr/share/mysql/japanese/errmsg.sys
usr/share/mysql/korean/
usr/share/mysql/korean/errmsg.sys
usr/share/mysql/magic
usr/share/mysql/maria_add_gis_sp.sql
usr/share/mysql/maria_add_gis_sp_bootstrap.sql
usr/share/mysql/mariadb.logrotate
usr/share/mysql/mariadb_performance_tables.sql
usr/share/mysql/mariadb_sys_schema.sql
usr/share/mysql/mariadb_system_tables.sql
usr/share/mysql/mariadb_system_tables_data.sql
usr/share/mysql/mariadb_test_data_timezone.sql
usr/share/mysql/mariadb_test_db.sql
usr/share/mysql/mini-benchmark
usr/share/mysql/mroonga/
usr/share/mysql/mroonga/AUTHORS
usr/share/mysql/mroonga/COPYING
usr/share/mysql/mroonga/install.sql
usr/share/mysql/mroonga/uninstall.sql
usr/share/mysql/norwegian-ny/
usr/share/mysql/norwegian-ny/errmsg.sys
usr/share/mysql/norwegian/
usr/share/mysql/norwegian/errmsg.sys
usr/share/mysql/policy/
usr/share/mysql/policy/apparmor/
usr/share/mysql/policy/apparmor/README
usr/share/mysql/policy/apparmor/usr.sbin.mysqld
usr/share/mysql/policy/apparmor/usr.sbin.mysqld.local
usr/share/mysql/policy/selinux/
usr/share/mysql/policy/selinux/README
usr/share/mysql/policy/selinux/mariadb-server.fc
usr/share/mysql/policy/selinux/mariadb-server.te
usr/share/mysql/policy/selinux/mariadb.te
usr/share/mysql/polish/
usr/share/mysql/polish/errmsg.sys
usr/share/mysql/portuguese/
usr/share/mysql/portuguese/errmsg.sys
usr/share/mysql/romanian/
usr/share/mysql/romanian/errmsg.sys
usr/share/mysql/russian/
usr/share/mysql/russian/errmsg.sys
usr/share/mysql/serbian/
usr/share/mysql/serbian/errmsg.sys
usr/share/mysql/slovak/
usr/share/mysql/slovak/errmsg.sys
usr/share/mysql/spanish/
usr/share/mysql/spanish/errmsg.sys
usr/share/mysql/swahili/
usr/share/mysql/swahili/errmsg.sys
usr/share/mysql/swedish/
usr/share/mysql/swedish/errmsg.sys
usr/share/mysql/ukrainian/
usr/share/mysql/ukrainian/errmsg.sys
usr/share/mysql/wsrep.cnf
usr/share/mysql/wsrep_notify
