# extracted/drone-runner-ssh-1.0.1-5/files
%FILES%
usr/
usr/bin/
usr/bin/drone-runner-ssh
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/drone-runner-ssh.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/drone-runner-ssh.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/drone-runner-ssh.conf
usr/share/
usr/share/doc/
usr/share/doc/drone-runner-ssh/
usr/share/doc/drone-runner-ssh/CHANGELOG.md
usr/share/doc/drone-runner-ssh/README.md
usr/share/licenses/
usr/share/licenses/drone-runner-ssh/
usr/share/licenses/drone-runner-ssh/LICENSE.md
usr/share/licenses/drone-runner-ssh/Polyform-Free-Trial.md
usr/share/licenses/drone-runner-ssh/Polyform-Noncommercial.md
usr/share/licenses/drone-runner-ssh/Polyform-Small-Business.md
