# extracted/netperf-2.7.0-6/files
%FILES%
usr/
usr/bin/
usr/bin/netperf
usr/bin/netserver
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/netserver.service
usr/share/
usr/share/info/
usr/share/info/netperf.info.gz
usr/share/licenses/
usr/share/licenses/netperf/
usr/share/licenses/netperf/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/netperf.1.gz
usr/share/man/man1/netserver.1.gz
