# extracted/drone-runner-docker-1.8.3-2/files
%FILES%
usr/
usr/bin/
usr/bin/drone-runner-docker
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/drone-runner-docker.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/drone-runner-docker.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/drone-runner-docker.conf
usr/share/
usr/share/doc/
usr/share/doc/drone-runner-docker/
usr/share/doc/drone-runner-docker/CHANGELOG.md
usr/share/doc/drone-runner-docker/README.md
usr/share/licenses/
usr/share/licenses/drone-runner-docker/
usr/share/licenses/drone-runner-docker/LICENSE.md
usr/share/licenses/drone-runner-docker/Polyform-Free-Trial.md
usr/share/licenses/drone-runner-docker/Polyform-Small-Business.md
