# extracted/slurm-llnl-24.11.0.1-1/files
%FILES%
etc/
etc/default/
etc/default/slurmctld
etc/default/slurmd
etc/default/slurmdbd
etc/slurm-llnl/
etc/slurm-llnl/cgroup.conf.example
etc/slurm-llnl/slurm.conf.example
etc/slurm-llnl/slurmdbd.conf.example
usr/
usr/bin/
usr/bin/sacct
usr/bin/sacctmgr
usr/bin/sackd
usr/bin/salloc
usr/bin/sattach
usr/bin/sbatch
usr/bin/sbcast
usr/bin/scancel
usr/bin/scontrol
usr/bin/scrontab
usr/bin/scrun
usr/bin/sdiag
usr/bin/sh5util
usr/bin/sinfo
usr/bin/slurmctld
usr/bin/slurmd
usr/bin/slurmdbd
usr/bin/slurmrestd
usr/bin/slurmstepd
usr/bin/sprio
usr/bin/squeue
usr/bin/sreport
usr/bin/srun
usr/bin/sshare
usr/bin/sstat
usr/bin/strigger
usr/bin/sview
usr/include/
usr/include/slurm/
usr/include/slurm/pmi.h
usr/include/slurm/slurm.h
usr/include/slurm/slurm_errno.h
usr/include/slurm/slurm_version.h
usr/include/slurm/slurmdb.h
usr/include/slurm/spank.h
usr/lib/
usr/lib/libslurm.so
usr/lib/libslurm.so.42
usr/lib/libslurm.so.42.0.0
usr/lib/slurm/
usr/lib/slurm/accounting_storage_ctld_relay.so
usr/lib/slurm/accounting_storage_mysql.so
usr/lib/slurm/accounting_storage_slurmdbd.so
usr/lib/slurm/acct_gather_energy_gpu.so
usr/lib/slurm/acct_gather_energy_ibmaem.so
usr/lib/slurm/acct_gather_energy_ipmi.so
usr/lib/slurm/acct_gather_energy_pm_counters.so
usr/lib/slurm/acct_gather_energy_rapl.so
usr/lib/slurm/acct_gather_energy_xcc.so
usr/lib/slurm/acct_gather_filesystem_lustre.so
usr/lib/slurm/acct_gather_interconnect_ofed.so
usr/lib/slurm/acct_gather_interconnect_sysfs.so
usr/lib/slurm/acct_gather_profile_hdf5.so
usr/lib/slurm/acct_gather_profile_influxdb.so
usr/lib/slurm/auth_munge.so
usr/lib/slurm/auth_none.so
usr/lib/slurm/burst_buffer_datawarp.so
usr/lib/slurm/burst_buffer_lua.so
usr/lib/slurm/certmgr_script.so
usr/lib/slurm/cgroup_v1.so
usr/lib/slurm/cgroup_v2.so
usr/lib/slurm/cli_filter_lua.so
usr/lib/slurm/cli_filter_syslog.so
usr/lib/slurm/cli_filter_user_defaults.so
usr/lib/slurm/cred_munge.so
usr/lib/slurm/cred_none.so
usr/lib/slurm/data_parser_v0_0_40.so
usr/lib/slurm/data_parser_v0_0_41.so
usr/lib/slurm/data_parser_v0_0_42.so
usr/lib/slurm/gpu_generic.so
usr/lib/slurm/gpu_nrt.so
usr/lib/slurm/gpu_nvidia.so
usr/lib/slurm/gpu_oneapi.so
usr/lib/slurm/gres_gpu.so
usr/lib/slurm/gres_mps.so
usr/lib/slurm/gres_nic.so
usr/lib/slurm/gres_shard.so
usr/lib/slurm/hash_k12.so
usr/lib/slurm/hash_sha3.so
usr/lib/slurm/job_container_tmpfs.so
usr/lib/slurm/job_submit_all_partitions.so
usr/lib/slurm/job_submit_defaults.so
usr/lib/slurm/job_submit_logging.so
usr/lib/slurm/job_submit_lua.so
usr/lib/slurm/job_submit_partition.so
usr/lib/slurm/job_submit_pbs.so
usr/lib/slurm/job_submit_require_timelimit.so
usr/lib/slurm/job_submit_throttle.so
usr/lib/slurm/jobacct_gather_cgroup.so
usr/lib/slurm/jobacct_gather_linux.so
usr/lib/slurm/jobcomp_elasticsearch.so
usr/lib/slurm/jobcomp_filetxt.so
usr/lib/slurm/jobcomp_lua.so
usr/lib/slurm/jobcomp_mysql.so
usr/lib/slurm/jobcomp_script.so
usr/lib/slurm/libslurm_pmi.so
usr/lib/slurm/libslurmfull.so
usr/lib/slurm/mcs_account.so
usr/lib/slurm/mcs_group.so
usr/lib/slurm/mcs_label.so
usr/lib/slurm/mcs_user.so
usr/lib/slurm/mpi_cray_shasta.so
usr/lib/slurm/mpi_pmi2.so
usr/lib/slurm/mpi_pmix.so
usr/lib/slurm/mpi_pmix_v5.so
usr/lib/slurm/node_features_helpers.so
usr/lib/slurm/node_features_knl_generic.so
usr/lib/slurm/openapi_slurmctld.so
usr/lib/slurm/openapi_slurmdbd.so
usr/lib/slurm/preempt_partition_prio.so
usr/lib/slurm/preempt_qos.so
usr/lib/slurm/prep_script.so
usr/lib/slurm/priority_basic.so
usr/lib/slurm/priority_multifactor.so
usr/lib/slurm/proctrack_cgroup.so
usr/lib/slurm/proctrack_linuxproc.so
usr/lib/slurm/proctrack_pgid.so
usr/lib/slurm/rest_auth_local.so
usr/lib/slurm/sched_backfill.so
usr/lib/slurm/sched_builtin.so
usr/lib/slurm/select_cons_tres.so
usr/lib/slurm/select_linear.so
usr/lib/slurm/serializer_json.so
usr/lib/slurm/serializer_url_encoded.so
usr/lib/slurm/serializer_yaml.so
usr/lib/slurm/site_factor_example.so
usr/lib/slurm/spank_pbs.so
usr/lib/slurm/src/
usr/lib/slurm/src/sattach/
usr/lib/slurm/src/sattach/sattach.wrapper.c
usr/lib/slurm/src/srun/
usr/lib/slurm/src/srun/srun.wrapper.c
usr/lib/slurm/switch_nvidia_imex.so
usr/lib/slurm/task_affinity.so
usr/lib/slurm/task_cgroup.so
usr/lib/slurm/tls_none.so
usr/lib/slurm/topology_3d_torus.so
usr/lib/slurm/topology_block.so
usr/lib/slurm/topology_default.so
usr/lib/slurm/topology_tree.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/slurmctld.service
usr/lib/systemd/system/slurmd.service
usr/lib/systemd/system/slurmdbd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/slurm-llnl.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/slurm-llnl.conf
usr/share/
usr/share/doc/
usr/share/doc/slurm-24.11.0/
usr/share/doc/slurm-24.11.0/html/
usr/share/doc/slurm-24.11.0/html/Slurm_Entity.pdf
usr/share/doc/slurm-24.11.0/html/Slurm_Individual.pdf
usr/share/doc/slurm-24.11.0/html/accounting.html
usr/share/doc/slurm-24.11.0/html/add.html
usr/share/doc/slurm-24.11.0/html/allocation_pies.gif
usr/share/doc/slurm-24.11.0/html/api.html
usr/share/doc/slurm-24.11.0/html/arch.gif
usr/share/doc/slurm-24.11.0/html/authentication.html
usr/share/doc/slurm-24.11.0/html/big_sys.html
usr/share/doc/slurm-24.11.0/html/burst_buffer.html
usr/share/doc/slurm-24.11.0/html/cg_hierarchy.jpg
usr/share/doc/slurm-24.11.0/html/cgroup_v2.html
usr/share/doc/slurm-24.11.0/html/cgroups.html
usr/share/doc/slurm-24.11.0/html/classic_fair_share.html
usr/share/doc/slurm-24.11.0/html/cli_filter_plugins.html
usr/share/doc/slurm-24.11.0/html/coding_style.pdf
usr/share/doc/slurm-24.11.0/html/configless_slurm.html
usr/share/doc/slurm-24.11.0/html/configurator.easy.html
usr/share/doc/slurm-24.11.0/html/configurator.html
usr/share/doc/slurm-24.11.0/html/cons_tres.html
usr/share/doc/slurm-24.11.0/html/cons_tres_share.html
usr/share/doc/slurm-24.11.0/html/containers.html
usr/share/doc/slurm-24.11.0/html/contributor.html
usr/share/doc/slurm-24.11.0/html/core_spec.html
usr/share/doc/slurm-24.11.0/html/cpu_management.html
usr/share/doc/slurm-24.11.0/html/dist_plane.html
usr/share/doc/slurm-24.11.0/html/documentation.html
usr/share/doc/slurm-24.11.0/html/dynamic_nodes.html
usr/share/doc/slurm-24.11.0/html/elasticsearch.html
usr/share/doc/slurm-24.11.0/html/entities.gif
usr/share/doc/slurm-24.11.0/html/example_usage.gif
usr/share/doc/slurm-24.11.0/html/extra_constraints.html
usr/share/doc/slurm-24.11.0/html/fair_tree.html
usr/share/doc/slurm-24.11.0/html/faq.html
usr/share/doc/slurm-24.11.0/html/federation.html
usr/share/doc/slurm-24.11.0/html/fonts.css
usr/share/doc/slurm-24.11.0/html/fonts.ttf
usr/share/doc/slurm-24.11.0/html/gang_scheduling.html
usr/share/doc/slurm-24.11.0/html/gres.html
usr/share/doc/slurm-24.11.0/html/gres_design.html
usr/share/doc/slurm-24.11.0/html/hdf5_job_outline.png
usr/share/doc/slurm-24.11.0/html/hdf5_profile_user_guide.html
usr/share/doc/slurm-24.11.0/html/hdf5_task_attr.png
usr/share/doc/slurm-24.11.0/html/heterogeneous_jobs.html
usr/share/doc/slurm-24.11.0/html/high_throughput.html
usr/share/doc/slurm-24.11.0/html/ibm_pe_fig1.png
usr/share/doc/slurm-24.11.0/html/ibm_pe_fig2.png
usr/share/doc/slurm-24.11.0/html/intel_knl.html
usr/share/doc/slurm-24.11.0/html/job_array.html
usr/share/doc/slurm-24.11.0/html/job_container_tmpfs.html
usr/share/doc/slurm-24.11.0/html/job_exit_code.html
usr/share/doc/slurm-24.11.0/html/job_launch.html
usr/share/doc/slurm-24.11.0/html/job_reason_codes.html
usr/share/doc/slurm-24.11.0/html/job_state_codes.html
usr/share/doc/slurm-24.11.0/html/job_submit_plugins.html
usr/share/doc/slurm-24.11.0/html/jobcomp_kafka.html
usr/share/doc/slurm-24.11.0/html/jobcomp_kafka_fig1.png
usr/share/doc/slurm-24.11.0/html/jquery.min.js
usr/share/doc/slurm-24.11.0/html/jwt.html
usr/share/doc/slurm-24.11.0/html/k_function.gif
usr/share/doc/slurm-24.11.0/html/kubernetes.html
usr/share/doc/slurm-24.11.0/html/licenses.html
usr/share/doc/slurm-24.11.0/html/mail.html
usr/share/doc/slurm-24.11.0/html/man_index.html
usr/share/doc/slurm-24.11.0/html/mc_support.gif
usr/share/doc/slurm-24.11.0/html/mc_support.html
usr/share/doc/slurm-24.11.0/html/mcs.html
usr/share/doc/slurm-24.11.0/html/mpi_guide.html
usr/share/doc/slurm-24.11.0/html/multi_cluster.html
usr/share/doc/slurm-24.11.0/html/network.html
usr/share/doc/slurm-24.11.0/html/network_failover.gif
usr/share/doc/slurm-24.11.0/html/network_federation.gif
usr/share/doc/slurm-24.11.0/html/network_multi_cluster.gif
usr/share/doc/slurm-24.11.0/html/network_srun.gif
usr/share/doc/slurm-24.11.0/html/network_standard.gif
usr/share/doc/slurm-24.11.0/html/node_lifecycle.png
usr/share/doc/slurm-24.11.0/html/nss_slurm.html
usr/share/doc/slurm-24.11.0/html/openapi.css
usr/share/doc/slurm-24.11.0/html/openapi_release_notes.html
usr/share/doc/slurm-24.11.0/html/overview.html
usr/share/doc/slurm-24.11.0/html/pam_slurm_adopt.html
usr/share/doc/slurm-24.11.0/html/plane_ex1.gif
usr/share/doc/slurm-24.11.0/html/plane_ex2.gif
usr/share/doc/slurm-24.11.0/html/plane_ex3.gif
usr/share/doc/slurm-24.11.0/html/plane_ex4.gif
usr/share/doc/slurm-24.11.0/html/plane_ex5.gif
usr/share/doc/slurm-24.11.0/html/plane_ex6.gif
usr/share/doc/slurm-24.11.0/html/plane_ex7.gif
usr/share/doc/slurm-24.11.0/html/platforms.html
usr/share/doc/slurm-24.11.0/html/plugins.html
usr/share/doc/slurm-24.11.0/html/power_save.html
usr/share/doc/slurm-24.11.0/html/preempt.html
usr/share/doc/slurm-24.11.0/html/prep_plugins.html
usr/share/doc/slurm-24.11.0/html/priority_multifactor.html
usr/share/doc/slurm-24.11.0/html/priority_multifactor3.html
usr/share/doc/slurm-24.11.0/html/programmer_guide.html
usr/share/doc/slurm-24.11.0/html/prolog_epilog.html
usr/share/doc/slurm-24.11.0/html/publications.html
usr/share/doc/slurm-24.11.0/html/qos.html
usr/share/doc/slurm-24.11.0/html/quickstart.html
usr/share/doc/slurm-24.11.0/html/quickstart_admin.html
usr/share/doc/slurm-24.11.0/html/related_software.html
usr/share/doc/slurm-24.11.0/html/release_notes.html
usr/share/doc/slurm-24.11.0/html/reservations.html
usr/share/doc/slurm-24.11.0/html/reset.css
usr/share/doc/slurm-24.11.0/html/resource_binding.html
usr/share/doc/slurm-24.11.0/html/resource_limits.html
usr/share/doc/slurm-24.11.0/html/rest.html
usr/share/doc/slurm-24.11.0/html/rest_api.html
usr/share/doc/slurm-24.11.0/html/rest_clients.html
usr/share/doc/slurm-24.11.0/html/rest_quickstart.html
usr/share/doc/slurm-24.11.0/html/rosetta.html
usr/share/doc/slurm-24.11.0/html/sched_config.html
usr/share/doc/slurm-24.11.0/html/schedmd.png
usr/share/doc/slurm-24.11.0/html/select_design.html
usr/share/doc/slurm-24.11.0/html/selinux.html
usr/share/doc/slurm-24.11.0/html/site_factor.html
usr/share/doc/slurm-24.11.0/html/slinky.html
usr/share/doc/slurm-24.11.0/html/slurm.css
usr/share/doc/slurm-24.11.0/html/slurm.html
usr/share/doc/slurm-24.11.0/html/slurm_logo.png
usr/share/doc/slurm-24.11.0/html/slurm_ug_agenda.html
usr/share/doc/slurm-24.11.0/html/squeue_color.png
usr/share/doc/slurm-24.11.0/html/style.css
usr/share/doc/slurm-24.11.0/html/topo_ex1.gif
usr/share/doc/slurm-24.11.0/html/topo_ex2.gif
usr/share/doc/slurm-24.11.0/html/topology.html
usr/share/doc/slurm-24.11.0/html/tres.html
usr/share/doc/slurm-24.11.0/html/troubleshoot.html
usr/share/doc/slurm-24.11.0/html/tutorial_intro_files.tar
usr/share/doc/slurm-24.11.0/html/upgrades.html
usr/share/doc/slurm-24.11.0/html/usage_pies.gif
usr/share/doc/slurm-24.11.0/html/user_permissions.html
usr/share/doc/slurm-24.11.0/html/wckey.html
usr/share/licenses/
usr/share/licenses/slurm-llnl/
usr/share/licenses/slurm-llnl/COPYING
usr/share/licenses/slurm-llnl/LICENSE.OpenSSL
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/sacct.1.gz
usr/share/man/man1/sacctmgr.1.gz
usr/share/man/man1/salloc.1.gz
usr/share/man/man1/sattach.1.gz
usr/share/man/man1/sbatch.1.gz
usr/share/man/man1/sbcast.1.gz
usr/share/man/man1/scancel.1.gz
usr/share/man/man1/scontrol.1.gz
usr/share/man/man1/scrontab.1.gz
usr/share/man/man1/scrun.1.gz
usr/share/man/man1/sdiag.1.gz
usr/share/man/man1/sh5util.1.gz
usr/share/man/man1/sinfo.1.gz
usr/share/man/man1/slurm.1.gz
usr/share/man/man1/sprio.1.gz
usr/share/man/man1/squeue.1.gz
usr/share/man/man1/sreport.1.gz
usr/share/man/man1/srun.1.gz
usr/share/man/man1/sshare.1.gz
usr/share/man/man1/sstat.1.gz
usr/share/man/man1/strigger.1.gz
usr/share/man/man1/sview.1.gz
usr/share/man/man5/
usr/share/man/man5/acct_gather.conf.5.gz
usr/share/man/man5/burst_buffer.conf.5.gz
usr/share/man/man5/cgroup.conf.5.gz
usr/share/man/man5/gres.conf.5.gz
usr/share/man/man5/helpers.conf.5.gz
usr/share/man/man5/job_container.conf.5.gz
usr/share/man/man5/knl.conf.5.gz
usr/share/man/man5/mpi.conf.5.gz
usr/share/man/man5/oci.conf.5.gz
usr/share/man/man5/slurm.conf.5.gz
usr/share/man/man5/slurmdbd.conf.5.gz
usr/share/man/man5/topology.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/sackd.8.gz
usr/share/man/man8/slurmctld.8.gz
usr/share/man/man8/slurmd.8.gz
usr/share/man/man8/slurmdbd.8.gz
usr/share/man/man8/slurmrestd.8.gz
usr/share/man/man8/slurmstepd.8.gz
usr/share/man/man8/spank.8.gz
var/
var/lib/
var/lib/slurm-llnl/
var/log/
var/log/slurm-llnl/
