# extracted/libratbag-0.18-1/files
%FILES%
usr/
usr/bin/
usr/bin/lur-command
usr/bin/ratbagctl
usr/bin/ratbagd
usr/include/
usr/include/liblur.h
usr/lib/
usr/lib/liblur.so
usr/lib/liblur.so.3
usr/lib/liblur.so.3.0.3
usr/lib/pkgconfig/
usr/lib/pkgconfig/liblur.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ratbagd.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.ratbag1.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.ratbag1.conf
usr/share/doc/
usr/share/doc/libratbag/
usr/share/doc/libratbag/html/
usr/share/doc/libratbag/html/.buildinfo
usr/share/doc/libratbag/html/_images/
usr/share/doc/libratbag/html/_images/logo.svg
usr/share/doc/libratbag/html/_sources/
usr/share/doc/libratbag/html/_sources/dbus.rst.txt
usr/share/doc/libratbag/html/_sources/index.rst.txt
usr/share/doc/libratbag/html/_static/
usr/share/doc/libratbag/html/_static/_sphinx_javascript_frameworks_compat.js
usr/share/doc/libratbag/html/_static/basic.css
usr/share/doc/libratbag/html/_static/css/
usr/share/doc/libratbag/html/_static/css/badge_only.css
usr/share/doc/libratbag/html/_static/css/fonts/
usr/share/doc/libratbag/html/_static/css/fonts/Roboto-Slab-Bold.woff
usr/share/doc/libratbag/html/_static/css/fonts/Roboto-Slab-Bold.woff2
usr/share/doc/libratbag/html/_static/css/fonts/Roboto-Slab-Regular.woff
usr/share/doc/libratbag/html/_static/css/fonts/Roboto-Slab-Regular.woff2
usr/share/doc/libratbag/html/_static/css/fonts/fontawesome-webfont.eot
usr/share/doc/libratbag/html/_static/css/fonts/fontawesome-webfont.svg
usr/share/doc/libratbag/html/_static/css/fonts/fontawesome-webfont.ttf
usr/share/doc/libratbag/html/_static/css/fonts/fontawesome-webfont.woff
usr/share/doc/libratbag/html/_static/css/fonts/fontawesome-webfont.woff2
usr/share/doc/libratbag/html/_static/css/fonts/lato-bold-italic.woff
usr/share/doc/libratbag/html/_static/css/fonts/lato-bold-italic.woff2
usr/share/doc/libratbag/html/_static/css/fonts/lato-bold.woff
usr/share/doc/libratbag/html/_static/css/fonts/lato-bold.woff2
usr/share/doc/libratbag/html/_static/css/fonts/lato-normal-italic.woff
usr/share/doc/libratbag/html/_static/css/fonts/lato-normal-italic.woff2
usr/share/doc/libratbag/html/_static/css/fonts/lato-normal.woff
usr/share/doc/libratbag/html/_static/css/fonts/lato-normal.woff2
usr/share/doc/libratbag/html/_static/css/theme.css
usr/share/doc/libratbag/html/_static/doctools.js
usr/share/doc/libratbag/html/_static/documentation_options.js
usr/share/doc/libratbag/html/_static/file.png
usr/share/doc/libratbag/html/_static/jquery.js
usr/share/doc/libratbag/html/_static/js/
usr/share/doc/libratbag/html/_static/js/badge_only.js
usr/share/doc/libratbag/html/_static/js/theme.js
usr/share/doc/libratbag/html/_static/language_data.js
usr/share/doc/libratbag/html/_static/minus.png
usr/share/doc/libratbag/html/_static/plus.png
usr/share/doc/libratbag/html/_static/pygments.css
usr/share/doc/libratbag/html/_static/searchtools.js
usr/share/doc/libratbag/html/_static/sphinx_highlight.js
usr/share/doc/libratbag/html/dbus.html
usr/share/doc/libratbag/html/genindex.html
usr/share/doc/libratbag/html/index.html
usr/share/doc/libratbag/html/objects.inv
usr/share/doc/libratbag/html/search.html
usr/share/doc/libratbag/html/searchindex.js
usr/share/libratbag/
usr/share/libratbag/asus-rog-buzzard.device
usr/share/libratbag/asus-rog-chakram-core.device
usr/share/libratbag/asus-rog-chakram-x.device
usr/share/libratbag/asus-rog-chakram.device
usr/share/libratbag/asus-rog-gladius2-origin-pink.device
usr/share/libratbag/asus-rog-gladius2-origin.device
usr/share/libratbag/asus-rog-gladius2.device
usr/share/libratbag/asus-rog-harpe-wireless.device
usr/share/libratbag/asus-rog-keris-wireless-aimpoint.device
usr/share/libratbag/asus-rog-keris-wireless.device
usr/share/libratbag/asus-rog-pugio.device
usr/share/libratbag/asus-rog-spatha-x.device
usr/share/libratbag/asus-rog-strix-carry.device
usr/share/libratbag/asus-rog-strix-impact.device
usr/share/libratbag/asus-rog-strix-impact2-wireless.device
usr/share/libratbag/asus-rog-strix-impact2.device
usr/share/libratbag/etekcity-scroll-alpha.device
usr/share/libratbag/glorious-model-d.device
usr/share/libratbag/glorious-model-o.device
usr/share/libratbag/gskill-MX-780.device
usr/share/libratbag/logitech-M325.device
usr/share/libratbag/logitech-M500s.device
usr/share/libratbag/logitech-M545.device
usr/share/libratbag/logitech-M570.device
usr/share/libratbag/logitech-M585-M590.device
usr/share/libratbag/logitech-M705.device
usr/share/libratbag/logitech-M720.device
usr/share/libratbag/logitech-MX-Anywhere-3.device
usr/share/libratbag/logitech-MX-Anywhere-3S.device
usr/share/libratbag/logitech-MX-Anywhere2.device
usr/share/libratbag/logitech-MX-Anywhere2S.device
usr/share/libratbag/logitech-MX-Anywhere3B.device
usr/share/libratbag/logitech-MX-Ergo.device
usr/share/libratbag/logitech-MX-Master-2S.device
usr/share/libratbag/logitech-MX-Master-3.device
usr/share/libratbag/logitech-MX-Master-3S-for-Business.device
usr/share/libratbag/logitech-MX-Master-3S.device
usr/share/libratbag/logitech-MX-Master.device
usr/share/libratbag/logitech-MX-Vertical.device
usr/share/libratbag/logitech-MX518.device
usr/share/libratbag/logitech-Marathon-M705.device
usr/share/libratbag/logitech-T650.device
usr/share/libratbag/logitech-Wireless-Touchpad.device
usr/share/libratbag/logitech-g-powerplay.device
usr/share/libratbag/logitech-g-pro-keyboard.device
usr/share/libratbag/logitech-g-pro-wireless.device
usr/share/libratbag/logitech-g-pro-x-wireless-superlight.device
usr/share/libratbag/logitech-g-pro.device
usr/share/libratbag/logitech-g102-g203.device
usr/share/libratbag/logitech-g300.device
usr/share/libratbag/logitech-g302.device
usr/share/libratbag/logitech-g303-shroud.device
usr/share/libratbag/logitech-g303.device
usr/share/libratbag/logitech-g305.device
usr/share/libratbag/logitech-g402.device
usr/share/libratbag/logitech-g403-hero.device
usr/share/libratbag/logitech-g403-wireless.device
usr/share/libratbag/logitech-g403.device
usr/share/libratbag/logitech-g413.device
usr/share/libratbag/logitech-g5-2007.device
usr/share/libratbag/logitech-g5.device
usr/share/libratbag/logitech-g500.device
usr/share/libratbag/logitech-g500s.device
usr/share/libratbag/logitech-g502-hero-wireless.device
usr/share/libratbag/logitech-g502-hero.device
usr/share/libratbag/logitech-g502-proteus-core.device
usr/share/libratbag/logitech-g502-proteus-spectrum.device
usr/share/libratbag/logitech-g502-x-wireless.device
usr/share/libratbag/logitech-g502-x.device
usr/share/libratbag/logitech-g513.device
usr/share/libratbag/logitech-g600.device
usr/share/libratbag/logitech-g602.device
usr/share/libratbag/logitech-g603.device
usr/share/libratbag/logitech-g604.device
usr/share/libratbag/logitech-g635.device
usr/share/libratbag/logitech-g7.device
usr/share/libratbag/logitech-g700-wireless.device
usr/share/libratbag/logitech-g700.device
usr/share/libratbag/logitech-g700s.device
usr/share/libratbag/logitech-g703-hero.device
usr/share/libratbag/logitech-g703.device
usr/share/libratbag/logitech-g733.device
usr/share/libratbag/logitech-g815.device
usr/share/libratbag/logitech-g9.device
usr/share/libratbag/logitech-g900.device
usr/share/libratbag/logitech-g903-hero.device
usr/share/libratbag/logitech-g903.device
usr/share/libratbag/logitech-g910.device
usr/share/libratbag/logitech-g915-tkl.device
usr/share/libratbag/logitech-g915.device
usr/share/libratbag/logitech-g935.device
usr/share/libratbag/logitech-g9x-Call-of-Duty-MW3-Edition.device
usr/share/libratbag/logitech-g9x-Original.device
usr/share/libratbag/marsgaming-mm4.device
usr/share/libratbag/nubwo-x7-spectrum.device
usr/share/libratbag/openinput.device
usr/share/libratbag/roccat-kone-emp.device
usr/share/libratbag/roccat-kone-pure.device
usr/share/libratbag/roccat-kone-xtd.device
usr/share/libratbag/sinowealth-0027.device
usr/share/libratbag/sinowealth-0029.device
usr/share/libratbag/sinowealth-0051.device
usr/share/libratbag/sinowealth-1007.device
usr/share/libratbag/steelseries-kinzu-v2-pro.device
usr/share/libratbag/steelseries-kinzu-v2.device
usr/share/libratbag/steelseries-kinzu-v3.device
usr/share/libratbag/steelseries-rival-100.device
usr/share/libratbag/steelseries-rival-310.device
usr/share/libratbag/steelseries-rival-600.device
usr/share/libratbag/steelseries-rival-650-wireless.device
usr/share/libratbag/steelseries-rival.device
usr/share/libratbag/steelseries-sensei-310.device
usr/share/libratbag/steelseries-sensei-raw.device
usr/share/licenses/
usr/share/licenses/libratbag/
usr/share/licenses/libratbag/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/lur-command.1.gz
usr/share/man/man1/ratbagctl.1.gz
usr/share/man/man8/
usr/share/man/man8/ratbagd.8.gz
