# extracted/archiso-82-1/files
%FILES%
usr/
usr/bin/
usr/bin/mkarchiso
usr/bin/run_archiso
usr/share/
usr/share/archiso/
usr/share/archiso/configs/
usr/share/archiso/configs/baseline/
usr/share/archiso/configs/baseline/airootfs/
usr/share/archiso/configs/baseline/airootfs/etc/
usr/share/archiso/configs/baseline/airootfs/etc/locale.conf
usr/share/archiso/configs/baseline/airootfs/etc/localtime
usr/share/archiso/configs/baseline/airootfs/etc/mkinitcpio.conf.d/
usr/share/archiso/configs/baseline/airootfs/etc/mkinitcpio.conf.d/archiso.conf
usr/share/archiso/configs/baseline/airootfs/etc/mkinitcpio.d/
usr/share/archiso/configs/baseline/airootfs/etc/mkinitcpio.d/linux.preset
usr/share/archiso/configs/baseline/airootfs/etc/shadow
usr/share/archiso/configs/baseline/airootfs/etc/ssh/
usr/share/archiso/configs/baseline/airootfs/etc/ssh/sshd_config.d/
usr/share/archiso/configs/baseline/airootfs/etc/ssh/sshd_config.d/10-archiso.conf
usr/share/archiso/configs/baseline/airootfs/etc/systemd/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/network.conf.d/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/network.conf.d/ipv6-privacy-extensions.conf
usr/share/archiso/configs/baseline/airootfs/etc/systemd/network/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/network/20-ethernet.network
usr/share/archiso/configs/baseline/airootfs/etc/systemd/resolved.conf.d/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/resolved.conf.d/archiso.conf
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system-generators/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system-generators/systemd-gpt-auto-generator
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-config.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-final.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-local.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-main.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-network.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/hv_fcopy_daemon.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/hv_kvp_daemon.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/hv_vss_daemon.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/sshd.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/systemd-networkd.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/systemd-resolved.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/vboxservice.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/vmtoolsd.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/multi-user.target.wants/vmware-vmblock-fuse.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/network-online.target.wants/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/network-online.target.wants/systemd-networkd-wait-online.service
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/sockets.target.wants/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/sockets.target.wants/systemd-networkd.socket
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/
usr/share/archiso/configs/baseline/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/wait-for-only-one-interface.conf
usr/share/archiso/configs/baseline/bootstrap_packages.x86_64
usr/share/archiso/configs/baseline/efiboot/
usr/share/archiso/configs/baseline/efiboot/loader/
usr/share/archiso/configs/baseline/efiboot/loader/entries/
usr/share/archiso/configs/baseline/efiboot/loader/entries/01-archiso-x86_64-linux.conf
usr/share/archiso/configs/baseline/efiboot/loader/loader.conf
usr/share/archiso/configs/baseline/grub/
usr/share/archiso/configs/baseline/grub/grub.cfg
usr/share/archiso/configs/baseline/grub/loopback.cfg
usr/share/archiso/configs/baseline/packages.x86_64
usr/share/archiso/configs/baseline/pacman.conf
usr/share/archiso/configs/baseline/profiledef.sh
usr/share/archiso/configs/baseline/syslinux/
usr/share/archiso/configs/baseline/syslinux/syslinux-linux.cfg
usr/share/archiso/configs/baseline/syslinux/syslinux.cfg
usr/share/archiso/configs/releng/
usr/share/archiso/configs/releng/airootfs/
usr/share/archiso/configs/releng/airootfs/etc/
usr/share/archiso/configs/releng/airootfs/etc/hostname
usr/share/archiso/configs/releng/airootfs/etc/locale.conf
usr/share/archiso/configs/releng/airootfs/etc/localtime
usr/share/archiso/configs/releng/airootfs/etc/mkinitcpio.conf.d/
usr/share/archiso/configs/releng/airootfs/etc/mkinitcpio.conf.d/archiso.conf
usr/share/archiso/configs/releng/airootfs/etc/mkinitcpio.d/
usr/share/archiso/configs/releng/airootfs/etc/mkinitcpio.d/linux.preset
usr/share/archiso/configs/releng/airootfs/etc/modprobe.d/
usr/share/archiso/configs/releng/airootfs/etc/modprobe.d/broadcom-wl.conf
usr/share/archiso/configs/releng/airootfs/etc/motd
usr/share/archiso/configs/releng/airootfs/etc/pacman.d/
usr/share/archiso/configs/releng/airootfs/etc/pacman.d/hooks/
usr/share/archiso/configs/releng/airootfs/etc/pacman.d/hooks/uncomment-mirrors.hook
usr/share/archiso/configs/releng/airootfs/etc/pacman.d/hooks/zzzz99-remove-custom-hooks-from-airootfs.hook
usr/share/archiso/configs/releng/airootfs/etc/passwd
usr/share/archiso/configs/releng/airootfs/etc/resolv.conf
usr/share/archiso/configs/releng/airootfs/etc/shadow
usr/share/archiso/configs/releng/airootfs/etc/ssh/
usr/share/archiso/configs/releng/airootfs/etc/ssh/sshd_config.d/
usr/share/archiso/configs/releng/airootfs/etc/ssh/sshd_config.d/10-archiso.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/
usr/share/archiso/configs/releng/airootfs/etc/systemd/journald.conf.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/journald.conf.d/volatile-storage.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/logind.conf.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/logind.conf.d/do-not-suspend.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/network.conf.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/network.conf.d/ipv6-privacy-extensions.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/network/
usr/share/archiso/configs/releng/airootfs/etc/systemd/network/20-ethernet.network
usr/share/archiso/configs/releng/airootfs/etc/systemd/network/20-wlan.network
usr/share/archiso/configs/releng/airootfs/etc/systemd/network/20-wwan.network
usr/share/archiso/configs/releng/airootfs/etc/systemd/resolved.conf.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/resolved.conf.d/archiso.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/system-generators/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system-generators/systemd-gpt-auto-generator
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/choose-mirror.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-config.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-final.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-local.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-main.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-network.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/dbus-org.freedesktop.ModemManager1.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/dbus-org.freedesktop.network1.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/dbus-org.freedesktop.resolve1.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/dbus-org.freedesktop.timesync1.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/etc-pacman.d-gnupg.mount
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/getty@tty1.service.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/livecd-alsa-unmuter.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/livecd-talk.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/ModemManager.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/choose-mirror.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/hv_fcopy_daemon.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/hv_kvp_daemon.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/hv_vss_daemon.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/iwd.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/livecd-talk.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/pacman-init.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/reflector.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/sshd.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/systemd-networkd.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/systemd-resolved.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/vboxservice.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/vmtoolsd.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/multi-user.target.wants/vmware-vmblock-fuse.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/network-online.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/network-online.target.wants/systemd-networkd-wait-online.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/pacman-init.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/reflector.service.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/reflector.service.d/archiso.conf
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sockets.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sockets.target.wants/pcscd.socket
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sockets.target.wants/systemd-networkd.socket
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sound.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sound.target.wants/livecd-alsa-unmuter.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sysinit.target.wants/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sysinit.target.wants/systemd-time-wait-sync.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/sysinit.target.wants/systemd-timesyncd.service
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/
usr/share/archiso/configs/releng/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/wait-for-only-one-interface.conf
usr/share/archiso/configs/releng/airootfs/etc/xdg/
usr/share/archiso/configs/releng/airootfs/etc/xdg/reflector/
usr/share/archiso/configs/releng/airootfs/etc/xdg/reflector/reflector.conf
usr/share/archiso/configs/releng/airootfs/root/
usr/share/archiso/configs/releng/airootfs/root/.automated_script.sh
usr/share/archiso/configs/releng/airootfs/root/.gnupg/
usr/share/archiso/configs/releng/airootfs/root/.gnupg/scdaemon.conf
usr/share/archiso/configs/releng/airootfs/root/.zlogin
usr/share/archiso/configs/releng/airootfs/usr/
usr/share/archiso/configs/releng/airootfs/usr/local/
usr/share/archiso/configs/releng/airootfs/usr/local/bin/
usr/share/archiso/configs/releng/airootfs/usr/local/bin/Installation_guide
usr/share/archiso/configs/releng/airootfs/usr/local/bin/choose-mirror
usr/share/archiso/configs/releng/airootfs/usr/local/bin/livecd-sound
usr/share/archiso/configs/releng/airootfs/usr/local/share/
usr/share/archiso/configs/releng/airootfs/usr/local/share/livecd-sound/
usr/share/archiso/configs/releng/airootfs/usr/local/share/livecd-sound/asound.conf.in
usr/share/archiso/configs/releng/bootstrap_packages.x86_64
usr/share/archiso/configs/releng/efiboot/
usr/share/archiso/configs/releng/efiboot/loader/
usr/share/archiso/configs/releng/efiboot/loader/entries/
usr/share/archiso/configs/releng/efiboot/loader/entries/01-archiso-x86_64-linux.conf
usr/share/archiso/configs/releng/efiboot/loader/entries/02-archiso-x86_64-speech-linux.conf
usr/share/archiso/configs/releng/efiboot/loader/entries/03-archiso-x86_64-memtest86+.conf
usr/share/archiso/configs/releng/efiboot/loader/loader.conf
usr/share/archiso/configs/releng/grub/
usr/share/archiso/configs/releng/grub/grub.cfg
usr/share/archiso/configs/releng/grub/loopback.cfg
usr/share/archiso/configs/releng/packages.x86_64
usr/share/archiso/configs/releng/pacman.conf
usr/share/archiso/configs/releng/profiledef.sh
usr/share/archiso/configs/releng/syslinux/
usr/share/archiso/configs/releng/syslinux/archiso_head.cfg
usr/share/archiso/configs/releng/syslinux/archiso_pxe-linux.cfg
usr/share/archiso/configs/releng/syslinux/archiso_pxe.cfg
usr/share/archiso/configs/releng/syslinux/archiso_sys-linux.cfg
usr/share/archiso/configs/releng/syslinux/archiso_sys.cfg
usr/share/archiso/configs/releng/syslinux/archiso_tail.cfg
usr/share/archiso/configs/releng/syslinux/splash.png
usr/share/archiso/configs/releng/syslinux/syslinux.cfg
usr/share/doc/
usr/share/doc/archiso/
usr/share/doc/archiso/AUTHORS.rst
usr/share/doc/archiso/CHANGELOG.rst
usr/share/doc/archiso/CONTRIBUTING.rst
usr/share/doc/archiso/README.profile.rst
usr/share/doc/archiso/README.rst
usr/share/doc/archiso/README.transfer.rst
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mkarchiso.1.gz
