# extracted/opensips-3.4.2-1/files
%FILES%
etc/
etc/opensips/
etc/opensips/opensips.cfg
etc/opensips/scenario_callcenter.xml
etc/opensips/tls/
etc/opensips/tls/README
etc/opensips/tls/ca.conf
etc/opensips/tls/request.conf
etc/opensips/tls/rootCA/
etc/opensips/tls/rootCA/cacert.pem
etc/opensips/tls/rootCA/certs/
etc/opensips/tls/rootCA/certs/01.pem
etc/opensips/tls/rootCA/index.txt
etc/opensips/tls/rootCA/private/
etc/opensips/tls/rootCA/private/cakey.pem
etc/opensips/tls/rootCA/serial
etc/opensips/tls/user.conf
etc/opensips/tls/user/
etc/opensips/tls/user/user-calist.pem
etc/opensips/tls/user/user-cert.pem
etc/opensips/tls/user/user-cert_req.pem
etc/opensips/tls/user/user-privkey.pem
usr/
usr/bin/
usr/bin/opensips
usr/bin/osipsconfig
usr/lib/
usr/lib/opensips/
usr/lib/opensips/modules/
usr/lib/opensips/modules/acc.so
usr/lib/opensips/modules/alias_db.so
usr/lib/opensips/modules/auth.so
usr/lib/opensips/modules/auth_aaa.so
usr/lib/opensips/modules/auth_db.so
usr/lib/opensips/modules/avpops.so
usr/lib/opensips/modules/b2b_entities.so
usr/lib/opensips/modules/b2b_logic.so
usr/lib/opensips/modules/b2b_sca.so
usr/lib/opensips/modules/b2b_sdp_demux.so
usr/lib/opensips/modules/benchmark.so
usr/lib/opensips/modules/cachedb_local.so
usr/lib/opensips/modules/cachedb_sql.so
usr/lib/opensips/modules/call_center.so
usr/lib/opensips/modules/call_control.so
usr/lib/opensips/modules/callops.so
usr/lib/opensips/modules/cfgutils.so
usr/lib/opensips/modules/clusterer.so
usr/lib/opensips/modules/db_cachedb.so
usr/lib/opensips/modules/db_flatstore.so
usr/lib/opensips/modules/db_http.so
usr/lib/opensips/modules/db_mysql.so
usr/lib/opensips/modules/db_postgres.so
usr/lib/opensips/modules/db_text.so
usr/lib/opensips/modules/db_unixodbc.so
usr/lib/opensips/modules/db_virtual.so
usr/lib/opensips/modules/dialog.so
usr/lib/opensips/modules/dispatcher.so
usr/lib/opensips/modules/diversion.so
usr/lib/opensips/modules/dns_cache.so
usr/lib/opensips/modules/domain.so
usr/lib/opensips/modules/domainpolicy.so
usr/lib/opensips/modules/drouting.so
usr/lib/opensips/modules/enum.so
usr/lib/opensips/modules/event_datagram.so
usr/lib/opensips/modules/event_flatstore.so
usr/lib/opensips/modules/event_route.so
usr/lib/opensips/modules/event_routing.so
usr/lib/opensips/modules/event_stream.so
usr/lib/opensips/modules/event_virtual.so
usr/lib/opensips/modules/event_xmlrpc.so
usr/lib/opensips/modules/exec.so
usr/lib/opensips/modules/fraud_detection.so
usr/lib/opensips/modules/freeswitch.so
usr/lib/opensips/modules/freeswitch_scripting.so
usr/lib/opensips/modules/gflags.so
usr/lib/opensips/modules/group.so
usr/lib/opensips/modules/h350.so
usr/lib/opensips/modules/httpd.so
usr/lib/opensips/modules/imc.so
usr/lib/opensips/modules/jsonrpc.so
usr/lib/opensips/modules/ldap.so
usr/lib/opensips/modules/load_balancer.so
usr/lib/opensips/modules/mangler.so
usr/lib/opensips/modules/mathops.so
usr/lib/opensips/modules/maxfwd.so
usr/lib/opensips/modules/media_exchange.so
usr/lib/opensips/modules/mediaproxy.so
usr/lib/opensips/modules/mi_datagram.so
usr/lib/opensips/modules/mi_fifo.so
usr/lib/opensips/modules/mi_html.so
usr/lib/opensips/modules/mi_http.so
usr/lib/opensips/modules/mi_script.so
usr/lib/opensips/modules/mid_registrar.so
usr/lib/opensips/modules/msilo.so
usr/lib/opensips/modules/msrp_gateway.so
usr/lib/opensips/modules/msrp_relay.so
usr/lib/opensips/modules/msrp_ua.so
usr/lib/opensips/modules/nat_traversal.so
usr/lib/opensips/modules/nathelper.so
usr/lib/opensips/modules/options.so
usr/lib/opensips/modules/path.so
usr/lib/opensips/modules/peering.so
usr/lib/opensips/modules/permissions.so
usr/lib/opensips/modules/pike.so
usr/lib/opensips/modules/presence.so
usr/lib/opensips/modules/presence_callinfo.so
usr/lib/opensips/modules/presence_xcapdiff.so
usr/lib/opensips/modules/presence_xml.so
usr/lib/opensips/modules/prometheus.so
usr/lib/opensips/modules/proto_bin.so
usr/lib/opensips/modules/proto_bins.so
usr/lib/opensips/modules/proto_hep.so
usr/lib/opensips/modules/proto_msrp.so
usr/lib/opensips/modules/proto_smpp.so
usr/lib/opensips/modules/proto_tls.so
usr/lib/opensips/modules/proto_ws.so
usr/lib/opensips/modules/qos.so
usr/lib/opensips/modules/qrouting.so
usr/lib/opensips/modules/rate_cacher.so
usr/lib/opensips/modules/ratelimit.so
usr/lib/opensips/modules/registrar.so
usr/lib/opensips/modules/rr.so
usr/lib/opensips/modules/rtp_relay.so
usr/lib/opensips/modules/rtpengine.so
usr/lib/opensips/modules/rtpproxy.so
usr/lib/opensips/modules/script_helper.so
usr/lib/opensips/modules/signaling.so
usr/lib/opensips/modules/sip_i.so
usr/lib/opensips/modules/sipcapture.so
usr/lib/opensips/modules/sipmsgops.so
usr/lib/opensips/modules/sl.so
usr/lib/opensips/modules/speeddial.so
usr/lib/opensips/modules/sql_cacher.so
usr/lib/opensips/modules/sst.so
usr/lib/opensips/modules/statistics.so
usr/lib/opensips/modules/status_report.so
usr/lib/opensips/modules/stun.so
usr/lib/opensips/modules/tcp_mgm.so
usr/lib/opensips/modules/textops.so
usr/lib/opensips/modules/tls_mgm.so
usr/lib/opensips/modules/tls_openssl.so
usr/lib/opensips/modules/tm.so
usr/lib/opensips/modules/topology_hiding.so
usr/lib/opensips/modules/tracer.so
usr/lib/opensips/modules/uac.so
usr/lib/opensips/modules/uac_auth.so
usr/lib/opensips/modules/uac_redirect.so
usr/lib/opensips/modules/uac_registrant.so
usr/lib/opensips/modules/userblacklist.so
usr/lib/opensips/modules/usrloc.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opensips.service
usr/share/
usr/share/doc/
usr/share/doc/opensips/
usr/share/doc/opensips/AUTHORS
usr/share/doc/opensips/INSTALL
usr/share/doc/opensips/NEWS
usr/share/doc/opensips/README
usr/share/doc/opensips/README-MODULES
usr/share/doc/opensips/README.acc
usr/share/doc/opensips/README.alias_db
usr/share/doc/opensips/README.auth
usr/share/doc/opensips/README.auth_aaa
usr/share/doc/opensips/README.auth_db
usr/share/doc/opensips/README.avpops
usr/share/doc/opensips/README.b2b_entities
usr/share/doc/opensips/README.b2b_logic
usr/share/doc/opensips/README.b2b_sca
usr/share/doc/opensips/README.b2b_sdp_demux
usr/share/doc/opensips/README.benchmark
usr/share/doc/opensips/README.cachedb_local
usr/share/doc/opensips/README.cachedb_sql
usr/share/doc/opensips/README.call_center
usr/share/doc/opensips/README.call_control
usr/share/doc/opensips/README.callops
usr/share/doc/opensips/README.cfgutils
usr/share/doc/opensips/README.clusterer
usr/share/doc/opensips/README.db_cachedb
usr/share/doc/opensips/README.db_flatstore
usr/share/doc/opensips/README.db_http
usr/share/doc/opensips/README.db_mysql
usr/share/doc/opensips/README.db_postgres
usr/share/doc/opensips/README.db_text
usr/share/doc/opensips/README.db_unixodbc
usr/share/doc/opensips/README.db_virtual
usr/share/doc/opensips/README.dialog
usr/share/doc/opensips/README.dispatcher
usr/share/doc/opensips/README.diversion
usr/share/doc/opensips/README.dns_cache
usr/share/doc/opensips/README.domain
usr/share/doc/opensips/README.domainpolicy
usr/share/doc/opensips/README.drouting
usr/share/doc/opensips/README.enum
usr/share/doc/opensips/README.event_datagram
usr/share/doc/opensips/README.event_flatstore
usr/share/doc/opensips/README.event_route
usr/share/doc/opensips/README.event_routing
usr/share/doc/opensips/README.event_stream
usr/share/doc/opensips/README.event_virtual
usr/share/doc/opensips/README.event_xmlrpc
usr/share/doc/opensips/README.exec
usr/share/doc/opensips/README.fraud_detection
usr/share/doc/opensips/README.freeswitch
usr/share/doc/opensips/README.freeswitch_scripting
usr/share/doc/opensips/README.gflags
usr/share/doc/opensips/README.group
usr/share/doc/opensips/README.h350
usr/share/doc/opensips/README.httpd
usr/share/doc/opensips/README.imc
usr/share/doc/opensips/README.jsonrpc
usr/share/doc/opensips/README.ldap
usr/share/doc/opensips/README.load_balancer
usr/share/doc/opensips/README.mangler
usr/share/doc/opensips/README.mathops
usr/share/doc/opensips/README.maxfwd
usr/share/doc/opensips/README.media_exchange
usr/share/doc/opensips/README.mediaproxy
usr/share/doc/opensips/README.mi_datagram
usr/share/doc/opensips/README.mi_fifo
usr/share/doc/opensips/README.mi_html
usr/share/doc/opensips/README.mi_http
usr/share/doc/opensips/README.mi_script
usr/share/doc/opensips/README.mid_registrar
usr/share/doc/opensips/README.msilo
usr/share/doc/opensips/README.msrp_gateway
usr/share/doc/opensips/README.msrp_relay
usr/share/doc/opensips/README.msrp_ua
usr/share/doc/opensips/README.nat_traversal
usr/share/doc/opensips/README.nathelper
usr/share/doc/opensips/README.options
usr/share/doc/opensips/README.path
usr/share/doc/opensips/README.peering
usr/share/doc/opensips/README.permissions
usr/share/doc/opensips/README.pike
usr/share/doc/opensips/README.presence
usr/share/doc/opensips/README.presence_callinfo
usr/share/doc/opensips/README.presence_xcapdiff
usr/share/doc/opensips/README.presence_xml
usr/share/doc/opensips/README.prometheus
usr/share/doc/opensips/README.proto_bin
usr/share/doc/opensips/README.proto_bins
usr/share/doc/opensips/README.proto_hep
usr/share/doc/opensips/README.proto_msrp
usr/share/doc/opensips/README.proto_smpp
usr/share/doc/opensips/README.proto_tls
usr/share/doc/opensips/README.proto_ws
usr/share/doc/opensips/README.qos
usr/share/doc/opensips/README.qrouting
usr/share/doc/opensips/README.rate_cacher
usr/share/doc/opensips/README.ratelimit
usr/share/doc/opensips/README.registrar
usr/share/doc/opensips/README.rr
usr/share/doc/opensips/README.rtp_relay
usr/share/doc/opensips/README.rtpengine
usr/share/doc/opensips/README.rtpproxy
usr/share/doc/opensips/README.script_helper
usr/share/doc/opensips/README.signaling
usr/share/doc/opensips/README.sip_i
usr/share/doc/opensips/README.sipcapture
usr/share/doc/opensips/README.sipmsgops
usr/share/doc/opensips/README.sl
usr/share/doc/opensips/README.speeddial
usr/share/doc/opensips/README.sql_cacher
usr/share/doc/opensips/README.sst
usr/share/doc/opensips/README.statistics
usr/share/doc/opensips/README.status_report
usr/share/doc/opensips/README.stun
usr/share/doc/opensips/README.tcp_mgm
usr/share/doc/opensips/README.textops
usr/share/doc/opensips/README.tls_mgm
usr/share/doc/opensips/README.tls_openssl
usr/share/doc/opensips/README.tm
usr/share/doc/opensips/README.topology_hiding
usr/share/doc/opensips/README.tracer
usr/share/doc/opensips/README.uac
usr/share/doc/opensips/README.uac_auth
usr/share/doc/opensips/README.uac_redirect
usr/share/doc/opensips/README.uac_registrant
usr/share/doc/opensips/README.userblacklist
usr/share/doc/opensips/README.usrloc
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/opensips.cfg.5.gz
usr/share/man/man8/
usr/share/man/man8/opensips.8.gz
usr/share/opensips/
usr/share/opensips/dbtext/
usr/share/opensips/dbtext/opensips/
usr/share/opensips/dbtext/opensips/acc
usr/share/opensips/dbtext/opensips/active_watchers
usr/share/opensips/dbtext/opensips/address
usr/share/opensips/dbtext/opensips/aliases
usr/share/opensips/dbtext/opensips/b2b_entities
usr/share/opensips/dbtext/opensips/b2b_logic
usr/share/opensips/dbtext/opensips/b2b_sca
usr/share/opensips/dbtext/opensips/cachedb
usr/share/opensips/dbtext/opensips/carrierfailureroute
usr/share/opensips/dbtext/opensips/carrierroute
usr/share/opensips/dbtext/opensips/cc_agents
usr/share/opensips/dbtext/opensips/cc_calls
usr/share/opensips/dbtext/opensips/cc_cdrs
usr/share/opensips/dbtext/opensips/cc_flows
usr/share/opensips/dbtext/opensips/closeddial
usr/share/opensips/dbtext/opensips/clusterer
usr/share/opensips/dbtext/opensips/cpl
usr/share/opensips/dbtext/opensips/dbaliases
usr/share/opensips/dbtext/opensips/dialog
usr/share/opensips/dbtext/opensips/dialplan
usr/share/opensips/dbtext/opensips/dispatcher
usr/share/opensips/dbtext/opensips/domain
usr/share/opensips/dbtext/opensips/domainpolicy
usr/share/opensips/dbtext/opensips/dr_carriers
usr/share/opensips/dbtext/opensips/dr_gateways
usr/share/opensips/dbtext/opensips/dr_groups
usr/share/opensips/dbtext/opensips/dr_partitions
usr/share/opensips/dbtext/opensips/dr_rules
usr/share/opensips/dbtext/opensips/emergency_report
usr/share/opensips/dbtext/opensips/emergency_routing
usr/share/opensips/dbtext/opensips/emergency_service_provider
usr/share/opensips/dbtext/opensips/fraud_detection
usr/share/opensips/dbtext/opensips/freeswitch
usr/share/opensips/dbtext/opensips/freeswitch_scripting
usr/share/opensips/dbtext/opensips/globalblacklist
usr/share/opensips/dbtext/opensips/grp
usr/share/opensips/dbtext/opensips/imc_members
usr/share/opensips/dbtext/opensips/imc_rooms
usr/share/opensips/dbtext/opensips/jwt_profiles
usr/share/opensips/dbtext/opensips/jwt_secrets
usr/share/opensips/dbtext/opensips/load_balancer
usr/share/opensips/dbtext/opensips/location
usr/share/opensips/dbtext/opensips/missed_calls
usr/share/opensips/dbtext/opensips/presentity
usr/share/opensips/dbtext/opensips/pua
usr/share/opensips/dbtext/opensips/qr_profiles
usr/share/opensips/dbtext/opensips/rc_clients
usr/share/opensips/dbtext/opensips/rc_demo_ratesheet
usr/share/opensips/dbtext/opensips/rc_ratesheets
usr/share/opensips/dbtext/opensips/rc_vendors
usr/share/opensips/dbtext/opensips/re_grp
usr/share/opensips/dbtext/opensips/registrant
usr/share/opensips/dbtext/opensips/rls_presentity
usr/share/opensips/dbtext/opensips/rls_watchers
usr/share/opensips/dbtext/opensips/route_tree
usr/share/opensips/dbtext/opensips/rtpengine
usr/share/opensips/dbtext/opensips/rtpproxy_sockets
usr/share/opensips/dbtext/opensips/silo
usr/share/opensips/dbtext/opensips/sip_trace
usr/share/opensips/dbtext/opensips/smpp
usr/share/opensips/dbtext/opensips/speed_dial
usr/share/opensips/dbtext/opensips/subscriber
usr/share/opensips/dbtext/opensips/tcp_mgm
usr/share/opensips/dbtext/opensips/tls_mgm
usr/share/opensips/dbtext/opensips/uri
usr/share/opensips/dbtext/opensips/userblacklist
usr/share/opensips/dbtext/opensips/usr_preferences
usr/share/opensips/dbtext/opensips/version
usr/share/opensips/dbtext/opensips/watchers
usr/share/opensips/dbtext/opensips/xcap
usr/share/opensips/menuconfig_templates/
usr/share/opensips/menuconfig_templates/opensips_loadbalancer.m4
usr/share/opensips/menuconfig_templates/opensips_loadbalancer_def.m4
usr/share/opensips/menuconfig_templates/opensips_residential.m4
usr/share/opensips/menuconfig_templates/opensips_residential_def.m4
usr/share/opensips/menuconfig_templates/opensips_trunking.m4
usr/share/opensips/menuconfig_templates/opensips_trunking_def.m4
usr/share/opensips/mysql/
usr/share/opensips/mysql/acc-create.sql
usr/share/opensips/mysql/alias_db-create.sql
usr/share/opensips/mysql/auth_db-create.sql
usr/share/opensips/mysql/authjwt-create.sql
usr/share/opensips/mysql/avpops-create.sql
usr/share/opensips/mysql/b2b-create.sql
usr/share/opensips/mysql/b2b_sca-create.sql
usr/share/opensips/mysql/cachedb_sql-create.sql
usr/share/opensips/mysql/call_center-create.sql
usr/share/opensips/mysql/carrierroute-create.sql
usr/share/opensips/mysql/closeddial-create.sql
usr/share/opensips/mysql/clusterer-create.sql
usr/share/opensips/mysql/cpl-create.sql
usr/share/opensips/mysql/db-migrate.sql
usr/share/opensips/mysql/dialog-create.sql
usr/share/opensips/mysql/dialplan-create.sql
usr/share/opensips/mysql/dispatcher-create.sql
usr/share/opensips/mysql/domain-create.sql
usr/share/opensips/mysql/domainpolicy-create.sql
usr/share/opensips/mysql/drouting-create.sql
usr/share/opensips/mysql/emergency-create.sql
usr/share/opensips/mysql/fraud_detection-create.sql
usr/share/opensips/mysql/freeswitch_scripting-create.sql
usr/share/opensips/mysql/group-create.sql
usr/share/opensips/mysql/imc-create.sql
usr/share/opensips/mysql/load_balancer-create.sql
usr/share/opensips/mysql/msilo-create.sql
usr/share/opensips/mysql/permissions-create.sql
usr/share/opensips/mysql/presence-create.sql
usr/share/opensips/mysql/qrouting-create.sql
usr/share/opensips/mysql/ratecacher-create.sql
usr/share/opensips/mysql/registrant-create.sql
usr/share/opensips/mysql/rls-create.sql
usr/share/opensips/mysql/rtpengine-create.sql
usr/share/opensips/mysql/rtpproxy-create.sql
usr/share/opensips/mysql/smpp-create.sql
usr/share/opensips/mysql/speeddial-create.sql
usr/share/opensips/mysql/standard-create.sql
usr/share/opensips/mysql/table-migrate.sql
usr/share/opensips/mysql/tcp_mgm-create.sql
usr/share/opensips/mysql/tls_mgm-create.sql
usr/share/opensips/mysql/tracer-create.sql
usr/share/opensips/mysql/userblacklist-create.sql
usr/share/opensips/mysql/usrloc-create.sql
usr/share/opensips/postgres/
usr/share/opensips/postgres/acc-create.sql
usr/share/opensips/postgres/alias_db-create.sql
usr/share/opensips/postgres/auth_db-create.sql
usr/share/opensips/postgres/authjwt-create.sql
usr/share/opensips/postgres/avpops-create.sql
usr/share/opensips/postgres/b2b-create.sql
usr/share/opensips/postgres/b2b_sca-create.sql
usr/share/opensips/postgres/cachedb_sql-create.sql
usr/share/opensips/postgres/call_center-create.sql
usr/share/opensips/postgres/carrierroute-create.sql
usr/share/opensips/postgres/closeddial-create.sql
usr/share/opensips/postgres/clusterer-create.sql
usr/share/opensips/postgres/cpl-create.sql
usr/share/opensips/postgres/dialog-create.sql
usr/share/opensips/postgres/dialplan-create.sql
usr/share/opensips/postgres/dispatcher-create.sql
usr/share/opensips/postgres/domain-create.sql
usr/share/opensips/postgres/domainpolicy-create.sql
usr/share/opensips/postgres/drouting-create.sql
usr/share/opensips/postgres/emergency-create.sql
usr/share/opensips/postgres/fraud_detection-create.sql
usr/share/opensips/postgres/freeswitch_scripting-create.sql
usr/share/opensips/postgres/group-create.sql
usr/share/opensips/postgres/imc-create.sql
usr/share/opensips/postgres/load_balancer-create.sql
usr/share/opensips/postgres/msilo-create.sql
usr/share/opensips/postgres/permissions-create.sql
usr/share/opensips/postgres/presence-create.sql
usr/share/opensips/postgres/qrouting-create.sql
usr/share/opensips/postgres/ratecacher-create.sql
usr/share/opensips/postgres/registrant-create.sql
usr/share/opensips/postgres/rls-create.sql
usr/share/opensips/postgres/rtpengine-create.sql
usr/share/opensips/postgres/rtpproxy-create.sql
usr/share/opensips/postgres/smpp-create.sql
usr/share/opensips/postgres/speeddial-create.sql
usr/share/opensips/postgres/standard-create.sql
usr/share/opensips/postgres/tcp_mgm-create.sql
usr/share/opensips/postgres/tls_mgm-create.sql
usr/share/opensips/postgres/tracer-create.sql
usr/share/opensips/postgres/userblacklist-create.sql
usr/share/opensips/postgres/usrloc-create.sql
