# extracted/forgejo-9.0.2-1/files
%FILES%
etc/
etc/forgejo/
etc/forgejo/app.ini
usr/
usr/bin/
usr/bin/forgejo
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/forgejo.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/forgejo.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/forgejo.conf
usr/share/
usr/share/licenses/
usr/share/licenses/forgejo/
usr/share/licenses/forgejo/LICENSE
