# extracted/babeld-1.13-1/files
%FILES%
etc/
etc/babeld.conf
usr/
usr/bin/
usr/bin/babeld
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/babeld.service
usr/share/
usr/share/licenses/
usr/share/licenses/babeld/
usr/share/licenses/babeld/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/babeld.8.gz
