# extracted/yubikey-full-disk-encryption-r155.dcaac00-2/files
%FILES%
etc/
etc/ykfde.conf
usr/
usr/bin/
usr/bin/ykfde-enroll
usr/bin/ykfde-format
usr/bin/ykfde-open
usr/lib/
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/ykfde
usr/lib/initcpio/install/
usr/lib/initcpio/install/ykfde
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ykfde-suspend.service
usr/lib/ykfde-suspend/
usr/lib/ykfde-suspend/initramfs-suspend
usr/lib/ykfde-suspend/ykfde-suspend
usr/share/
usr/share/doc/
usr/share/doc/ykfde/
usr/share/doc/ykfde/README.md
