# extracted/matterbridge-1.26.0-2/files
%FILES%
etc/
etc/matterbridge.toml
usr/
usr/bin/
usr/bin/matterbridge
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/matterbridge.service
