# extracted/php-legacy-fpm-8.2.26-1/files
%FILES%
etc/
etc/php-legacy/
etc/php-legacy/php-fpm.conf
etc/php-legacy/php-fpm.d/
etc/php-legacy/php-fpm.d/www.conf
usr/
usr/bin/
usr/bin/php-fpm-legacy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/php-fpm-legacy.service
usr/lib/systemd/system/php-legacy-fpm.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/php-fpm-legacy.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/php-fpm-legacy.8.gz
usr/share/php-legacy/
usr/share/php-legacy/fpm/
usr/share/php-legacy/fpm/status.html
