# extracted/toxcore-1:0.2.19-1/files
%FILES%
etc/
etc/tox-bootstrapd.conf
usr/
usr/bin/
usr/bin/DHT_bootstrap
usr/bin/tox-bootstrapd
usr/include/
usr/include/tox/
usr/include/tox/tox.h
usr/include/tox/tox_dispatch.h
usr/include/tox/tox_events.h
usr/include/tox/toxav.h
usr/include/tox/toxencryptsave.h
usr/lib/
usr/lib/libtoxcore.so
usr/lib/libtoxcore.so.2
usr/lib/libtoxcore.so.2.19.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/toxcore.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tox-bootstrapd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/toxcore.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/toxcore.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/tox-bootstrapd
