# extracted/lemurs-0.3.2-3/files
%FILES%
etc/
etc/lemurs/
etc/lemurs/config.toml
etc/lemurs/wayland/
etc/lemurs/wms/
etc/lemurs/xsetup.sh
etc/pam.d/
etc/pam.d/lemurs
usr/
usr/bin/
usr/bin/lemurs
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lemurs.service
usr/share/
usr/share/doc/
usr/share/doc/lemurs/
usr/share/doc/lemurs/README.md
usr/share/licenses/
usr/share/licenses/lemurs/
usr/share/licenses/lemurs/MIT
