# extracted/miniflux-2.2.3-1/files
%FILES%
etc/
etc/miniflux.conf
usr/
usr/bin/
usr/bin/miniflux
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/miniflux.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/miniflux.conf
usr/share/
usr/share/doc/
usr/share/doc/miniflux/
usr/share/doc/miniflux/ChangeLog
usr/share/doc/miniflux/README.md
usr/share/doc/miniflux/SECURITY.md
usr/share/licenses/
usr/share/licenses/miniflux/
usr/share/licenses/miniflux/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/miniflux.1.gz
