# extracted/vault-1.18.2-1/files
%FILES%
etc/
etc/default/
etc/default/vault
etc/vault.hcl
usr/
usr/bin/
usr/bin/vault
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vault.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/vault.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/vault.conf
usr/share/
usr/share/doc/
usr/share/doc/vault/
usr/share/doc/vault/CHANGELOG.md
usr/share/doc/vault/CONTRIBUTING.md
usr/share/doc/vault/README.md
usr/share/licenses/
usr/share/licenses/vault/
usr/share/licenses/vault/LICENSE
