# extracted/uptimed-0.4.6-4/files
%FILES%
etc/
etc/uptimed.conf
usr/
usr/bin/
usr/bin/uprecords
usr/bin/uptimed
usr/lib/
usr/lib/libuptimed.so
usr/lib/libuptimed.so.0
usr/lib/libuptimed.so.0.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/uptimed.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/uprecords.1.gz
usr/share/man/man8/
usr/share/man/man8/uptimed.8.gz
var/
var/spool/
var/spool/uptimed/
