# extracted/opendmarc-1.4.2-4/files
%FILES%
etc/
etc/opendmarc/
etc/opendmarc/opendmarc.conf
usr/
usr/bin/
usr/bin/opendmarc
usr/bin/opendmarc-check
usr/bin/opendmarc-expire
usr/bin/opendmarc-import
usr/bin/opendmarc-importstats
usr/bin/opendmarc-params
usr/bin/opendmarc-reports
usr/include/
usr/include/opendmarc/
usr/include/opendmarc/dmarc.h
usr/lib/
usr/lib/libopendmarc.so
usr/lib/libopendmarc.so.2
usr/lib/libopendmarc.so.2.0.3
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opendmarc.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/opendmarc.conf
usr/share/
usr/share/doc/
usr/share/doc/opendmarc/
usr/share/doc/opendmarc/LICENSE.Sendmail
usr/share/doc/opendmarc/README
usr/share/doc/opendmarc/README.rddmarc
usr/share/doc/opendmarc/README.schema
usr/share/doc/opendmarc/dmarc_policy_t.html
usr/share/doc/opendmarc/dmarcfail.py
usr/share/doc/opendmarc/index.html
usr/share/doc/opendmarc/mkdmarc
usr/share/doc/opendmarc/mysql_ip6.c
usr/share/doc/opendmarc/opendmarc
usr/share/doc/opendmarc/opendmarc.conf.sample
usr/share/doc/opendmarc/opendmarc.service.in
usr/share/doc/opendmarc/opendmarc.spec.in
usr/share/doc/opendmarc/opendmarc_dns_fake_record.html
usr/share/doc/opendmarc/opendmarc_get_policy_to_enforce.html
usr/share/doc/opendmarc/opendmarc_lib_t.html
usr/share/doc/opendmarc/opendmarc_policy_connect_clear.html
usr/share/doc/opendmarc/opendmarc_policy_connect_init.html
usr/share/doc/opendmarc/opendmarc_policy_connect_rset.html
usr/share/doc/opendmarc/opendmarc_policy_connect_shutdown.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_adkim.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_alignment.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_aspf.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_fo.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_p.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_pct.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_rf.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_rua.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_ruf.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_sp.html
usr/share/doc/opendmarc/opendmarc_policy_fetch_utilized_domain.html
usr/share/doc/opendmarc/opendmarc_policy_library_init.html
usr/share/doc/opendmarc/opendmarc_policy_library_shutdown.html
usr/share/doc/opendmarc/opendmarc_policy_parse_dmarc.html
usr/share/doc/opendmarc/opendmarc_policy_query_dmarc.html
usr/share/doc/opendmarc/opendmarc_policy_status_to_str.html
usr/share/doc/opendmarc/opendmarc_policy_store_dkim.html
usr/share/doc/opendmarc/opendmarc_policy_store_dmarc.html
usr/share/doc/opendmarc/opendmarc_policy_store_from_domain.html
usr/share/doc/opendmarc/opendmarc_policy_store_spf.html
usr/share/doc/opendmarc/opendmarc_policy_to_buf.html
usr/share/doc/opendmarc/opendmarc_spf_test.html
usr/share/doc/opendmarc/opendmarc_status_t.html
usr/share/doc/opendmarc/opendmarc_tld_read_file.html
usr/share/doc/opendmarc/opendmarc_tld_shutdown.html
usr/share/doc/opendmarc/opendmarc_util_clearargv.html
usr/share/doc/opendmarc/opendmarc_xml.html
usr/share/doc/opendmarc/opendmarc_xml_parse.html
usr/share/doc/opendmarc/overview.html
usr/share/doc/opendmarc/rddmarc
usr/share/doc/opendmarc/schema.mysql
usr/share/licenses/
usr/share/licenses/opendmarc/
usr/share/licenses/opendmarc/LICENSE
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/opendmarc.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/opendmarc-check.8.gz
usr/share/man/man8/opendmarc-expire.8.gz
usr/share/man/man8/opendmarc-import.8.gz
usr/share/man/man8/opendmarc-importstats.8.gz
usr/share/man/man8/opendmarc-params.8.gz
usr/share/man/man8/opendmarc-reports.8.gz
usr/share/man/man8/opendmarc.8.gz
