# extracted/redis-7.2.5-1/files
%FILES%
etc/
etc/redis/
etc/redis/redis.conf
etc/redis/sentinel.conf
usr/
usr/bin/
usr/bin/redis-benchmark
usr/bin/redis-check-aof
usr/bin/redis-check-rdb
usr/bin/redis-cli
usr/bin/redis-sentinel
usr/bin/redis-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/redis-sentinel.service
usr/lib/systemd/system/redis.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/redis.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/redis.conf
usr/share/
usr/share/licenses/
usr/share/licenses/redis/
usr/share/licenses/redis/LICENSE
