# extracted/powerdns-4.9.2-1/files
%FILES%
etc/
etc/powerdns/
etc/powerdns/ixfrdist.example.yml
etc/powerdns/pdns.conf
usr/
usr/bin/
usr/bin/calidns
usr/bin/dnsbulktest
usr/bin/dnsgram
usr/bin/dnspcap2calidns
usr/bin/dnspcap2protobuf
usr/bin/dnsreplay
usr/bin/dnsscan
usr/bin/dnsscope
usr/bin/dnstcpbench
usr/bin/dnswasher
usr/bin/dumresp
usr/bin/ixfrdist
usr/bin/ixplore
usr/bin/nproxy
usr/bin/nsec3dig
usr/bin/pdns_control
usr/bin/pdns_notify
usr/bin/pdns_server
usr/bin/pdnsutil
usr/bin/saxfr
usr/bin/sdig
usr/bin/stubquery
usr/bin/zone2json
usr/bin/zone2ldap
usr/bin/zone2sql
usr/lib/
usr/lib/powerdns/
usr/lib/powerdns/libbindbackend.so
usr/lib/powerdns/libgeoipbackend.so
usr/lib/powerdns/libgmysqlbackend.so
usr/lib/powerdns/libgodbcbackend.so
usr/lib/powerdns/libgpgsqlbackend.so
usr/lib/powerdns/libgsqlite3backend.so
usr/lib/powerdns/libldapbackend.so
usr/lib/powerdns/liblmdbbackend.so
usr/lib/powerdns/liblua2backend.so
usr/lib/powerdns/libpipebackend.so
usr/lib/powerdns/libremotebackend.so
usr/lib/powerdns/libtinydnsbackend.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ixfrdist.service
usr/lib/systemd/system/ixfrdist@.service
usr/lib/systemd/system/pdns.service
usr/lib/systemd/system/pdns@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/powerdns.conf
usr/share/
usr/share/doc/
usr/share/doc/powerdns/
usr/share/doc/powerdns/3.4.0_to_4.0.0_schema.sqlite3.sql
usr/share/doc/powerdns/3.4.0_to_4.1.0_schema.mysql.sql
usr/share/doc/powerdns/3.4.0_to_4.1.0_schema.pgsql.sql
usr/share/doc/powerdns/4.0.0_to_4.2.0_schema.mssql.sql
usr/share/doc/powerdns/4.0.0_to_4.2.0_schema.sqlite3.sql
usr/share/doc/powerdns/4.1.0_to_4.2.0_schema.mysql.sql
usr/share/doc/powerdns/4.1.0_to_4.2.0_schema.pgsql.sql
usr/share/doc/powerdns/4.2.0_to_4.3.0_schema.mssql.sql
usr/share/doc/powerdns/4.2.0_to_4.3.0_schema.mysql.sql
usr/share/doc/powerdns/4.2.0_to_4.3.0_schema.pgsql.sql
usr/share/doc/powerdns/4.2.0_to_4.3.0_schema.sqlite3.sql
usr/share/doc/powerdns/4.3.0_to_4.3.1_schema.sqlite3.sql
usr/share/doc/powerdns/4.3.0_to_4.7.0_schema.mssql.sql
usr/share/doc/powerdns/4.3.0_to_4.7.0_schema.mysql.sql
usr/share/doc/powerdns/4.3.0_to_4.7.0_schema.pgsql.sql
usr/share/doc/powerdns/4.3.1_to_4.7.0_schema.sqlite3.sql
usr/share/doc/powerdns/bind-dnssec.4.2.0_to_4.3.0_schema.sqlite3.sql
usr/share/doc/powerdns/bind-dnssec.schema.sqlite3.sql
usr/share/doc/powerdns/dnsdomain2.schema
usr/share/doc/powerdns/dnssec-3.x_to_3.4.0_schema.mysql.sql
usr/share/doc/powerdns/dnssec-3.x_to_3.4.0_schema.pgsql.sql
usr/share/doc/powerdns/dnssec-3.x_to_3.4.0_schema.sqlite3.sql
usr/share/doc/powerdns/enable-foreign-keys.mysql.sql
usr/share/doc/powerdns/nodnssec-3.x_to_3.4.0_schema.mysql.sql
usr/share/doc/powerdns/nodnssec-3.x_to_3.4.0_schema.pgsql.sql
usr/share/doc/powerdns/nodnssec-3.x_to_3.4.0_schema.sqlite3.sql
usr/share/doc/powerdns/pdns-domaininfo.schema
usr/share/doc/powerdns/schema.mssql.sql
usr/share/doc/powerdns/schema.mysql.sql
usr/share/doc/powerdns/schema.pgsql.sql
usr/share/doc/powerdns/schema.sqlite3.sql
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/calidns.1.gz
usr/share/man/man1/dnsbulktest.1.gz
usr/share/man/man1/dnsgram.1.gz
usr/share/man/man1/dnspcap2calidns.1.gz
usr/share/man/man1/dnspcap2protobuf.1.gz
usr/share/man/man1/dnsreplay.1.gz
usr/share/man/man1/dnsscan.1.gz
usr/share/man/man1/dnsscope.1.gz
usr/share/man/man1/dnstcpbench.1.gz
usr/share/man/man1/dnswasher.1.gz
usr/share/man/man1/dumresp.1.gz
usr/share/man/man1/ixfrdist.1.gz
usr/share/man/man1/ixplore.1.gz
usr/share/man/man1/nproxy.1.gz
usr/share/man/man1/nsec3dig.1.gz
usr/share/man/man1/pdns_control.1.gz
usr/share/man/man1/pdns_notify.1.gz
usr/share/man/man1/pdns_server.1.gz
usr/share/man/man1/pdnsutil.1.gz
usr/share/man/man1/saxfr.1.gz
usr/share/man/man1/sdig.1.gz
usr/share/man/man1/zone2json.1.gz
usr/share/man/man1/zone2ldap.1.gz
usr/share/man/man1/zone2sql.1.gz
usr/share/man/man5/
usr/share/man/man5/ixfrdist.yml.5.gz
