# extracted/miniserve-0.28.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/miniserve
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/miniserve@.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/miniserve
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/miniserve.fish
usr/share/licenses/
usr/share/licenses/miniserve/
usr/share/licenses/miniserve/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/miniserve.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_miniserve
