# extracted/monero-0.18.3.4-2/files
%FILES%
etc/
etc/monerod.conf
usr/
usr/bin/
usr/bin/monero-blockchain-ancestry
usr/bin/monero-blockchain-depth
usr/bin/monero-blockchain-export
usr/bin/monero-blockchain-import
usr/bin/monero-blockchain-mark-spent-outputs
usr/bin/monero-blockchain-prune
usr/bin/monero-blockchain-prune-known-spent-data
usr/bin/monero-blockchain-stats
usr/bin/monero-blockchain-usage
usr/bin/monero-gen-ssl-cert
usr/bin/monero-gen-trusted-multisig
usr/bin/monero-wallet-cli
usr/bin/monero-wallet-rpc
usr/bin/monerod
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/monerod.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/monero.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/monero.conf
usr/share/
usr/share/licenses/
usr/share/licenses/monero/
usr/share/licenses/monero/LICENSE
