# extracted/yrd-0.5.3-8/files
%FILES%
usr/
usr/bin/
usr/bin/yrd
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/PKG-INFO
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/SOURCES.txt
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/dependency_links.txt
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/entry_points.txt
usr/lib/python3.12/site-packages/yrd-0.5.3-py3.12.egg-info/top_level.txt
usr/lib/python3.12/site-packages/yrd/
usr/lib/python3.12/site-packages/yrd/__init__.py
usr/lib/python3.12/site-packages/yrd/__pycache__/
usr/lib/python3.12/site-packages/yrd/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/arg.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/arg.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/const.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/const.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/core.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/core.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/peer.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/peer.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/start.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/start.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/tests.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/tests.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/top.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/top.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/utils.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/utils.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/xcjdns.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/__pycache__/xcjdns.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/arg.py
usr/lib/python3.12/site-packages/yrd/cjdns/
usr/lib/python3.12/site-packages/yrd/cjdns/__init__.py
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/bencode.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/bencode.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/cjdns.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/__pycache__/cjdns.cpython-312.pyc
usr/lib/python3.12/site-packages/yrd/cjdns/bencode.py
usr/lib/python3.12/site-packages/yrd/cjdns/cjdns.py
usr/lib/python3.12/site-packages/yrd/const.py
usr/lib/python3.12/site-packages/yrd/core.py
usr/lib/python3.12/site-packages/yrd/peer.py
usr/lib/python3.12/site-packages/yrd/start.py
usr/lib/python3.12/site-packages/yrd/tests.py
usr/lib/python3.12/site-packages/yrd/top.py
usr/lib/python3.12/site-packages/yrd/utils.py
usr/lib/python3.12/site-packages/yrd/xcjdns.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/yrd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/yrd.conf
usr/share/
usr/share/doc/
usr/share/doc/yrd/
usr/share/doc/yrd/README.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/yrd.1.gz
