# extracted/piping-server-0.18.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/piping-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/piping-server.service
usr/share/
usr/share/doc/
usr/share/doc/piping-server/
usr/share/doc/piping-server/README.md
usr/share/licenses/
usr/share/licenses/piping-server/
usr/share/licenses/piping-server/LICENSE
