# extracted/mailman-web-0.0.9-4/files
%FILES%
etc/
etc/uwsgi/
etc/uwsgi/mailman-web.ini
etc/webapps/
etc/webapps/mailman-web/
etc/webapps/mailman-web/settings.py
etc/webapps/mailman-web/urls.py
usr/
usr/bin/
usr/bin/mailman-web
usr/lib/
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/METADATA
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/RECORD
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/WHEEL
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/entry_points.txt
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/licenses/
usr/lib/python3.13/site-packages/mailman_web-0.0.9.dist-info/licenses/LICENSE.txt
usr/lib/python3.13/site-packages/mailman_web/
usr/lib/python3.13/site-packages/mailman_web/__init__.py
usr/lib/python3.13/site-packages/mailman_web/__pycache__/
usr/lib/python3.13/site-packages/mailman_web/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/manage.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/manage.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/urls.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/urls.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/wsgi.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/__pycache__/wsgi.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/manage.py
usr/lib/python3.13/site-packages/mailman_web/manage.py.orig
usr/lib/python3.13/site-packages/mailman_web/settings/
usr/lib/python3.13/site-packages/mailman_web/settings/__init__.py
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/base.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/base.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/dev.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/dev.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/mailman.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/__pycache__/mailman.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/settings/base.py
usr/lib/python3.13/site-packages/mailman_web/settings/base.py.orig
usr/lib/python3.13/site-packages/mailman_web/settings/dev.py
usr/lib/python3.13/site-packages/mailman_web/settings/mailman.py
usr/lib/python3.13/site-packages/mailman_web/tests/
usr/lib/python3.13/site-packages/mailman_web/tests/__init__.py
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/settings.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/settings.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/test_basic.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/__pycache__/test_basic.cpython-313.pyc
usr/lib/python3.13/site-packages/mailman_web/tests/settings.py
usr/lib/python3.13/site-packages/mailman_web/tests/test_basic.py
usr/lib/python3.13/site-packages/mailman_web/urls.py
usr/lib/python3.13/site-packages/mailman_web/wsgi.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mailman-web-qcluster.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mailman-web.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mailman-web.conf
usr/share/
usr/share/doc/
usr/share/doc/mailman-web/
usr/share/doc/mailman-web/README.rst
var/
var/lib/
var/lib/mailman-web/
var/lib/mailman-web/data/
var/lib/mailman-web/static/
var/log/
var/log/mailman-web/
