# extracted/backuppc-4.4.0-8/files
%FILES%
etc/
etc/backuppc/
etc/backuppc/backuppc.users
etc/backuppc/config.pl
etc/backuppc/config.pl.sample
etc/backuppc/hosts
etc/backuppc/hosts.sample
etc/httpd/
etc/httpd/conf/
etc/httpd/conf/extra/
etc/httpd/conf/extra/backuppc.conf
etc/profile.d/
etc/profile.d/backuppc.csh
etc/profile.d/backuppc.sh
usr/
usr/lib/
usr/lib/perl5/
usr/lib/perl5/5.40/
usr/lib/perl5/5.40/vendor_perl/
usr/lib/perl5/5.40/vendor_perl/BackupPC/
usr/lib/perl5/5.40/vendor_perl/BackupPC/XS.pm
usr/lib/perl5/5.40/vendor_perl/auto/
usr/lib/perl5/5.40/vendor_perl/auto/BackupPC/
usr/lib/perl5/5.40/vendor_perl/auto/BackupPC/XS/
usr/lib/perl5/5.40/vendor_perl/auto/BackupPC/XS/XS.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/backuppc.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/backuppc.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/backuppc.conf
usr/share/
usr/share/backuppc/
usr/share/backuppc/bin/
usr/share/backuppc/bin/BackupPC
usr/share/backuppc/bin/BackupPC_Admin_SCGI
usr/share/backuppc/bin/BackupPC_archive
usr/share/backuppc/bin/BackupPC_archiveHost
usr/share/backuppc/bin/BackupPC_archiveStart
usr/share/backuppc/bin/BackupPC_attribPrint
usr/share/backuppc/bin/BackupPC_backupDelete
usr/share/backuppc/bin/BackupPC_backupDuplicate
usr/share/backuppc/bin/BackupPC_dump
usr/share/backuppc/bin/BackupPC_fixupBackupSummary
usr/share/backuppc/bin/BackupPC_fsck
usr/share/backuppc/bin/BackupPC_ls
usr/share/backuppc/bin/BackupPC_migrateV3toV4
usr/share/backuppc/bin/BackupPC_nightly
usr/share/backuppc/bin/BackupPC_poolCntPrint
usr/share/backuppc/bin/BackupPC_refCountUpdate
usr/share/backuppc/bin/BackupPC_restore
usr/share/backuppc/bin/BackupPC_rrdUpdate
usr/share/backuppc/bin/BackupPC_sendEmail
usr/share/backuppc/bin/BackupPC_serverMesg
usr/share/backuppc/bin/BackupPC_tarCreate
usr/share/backuppc/bin/BackupPC_tarExtract
usr/share/backuppc/bin/BackupPC_zcat
usr/share/backuppc/bin/BackupPC_zipCreate
usr/share/backuppc/bin/rsync_bpc
usr/share/backuppc/cgi-bin/
usr/share/backuppc/cgi-bin/BackupPC_Admin
usr/share/backuppc/html/
usr/share/backuppc/html/0000000.gif
usr/share/backuppc/html/0000011.gif
usr/share/backuppc/html/0001000.gif
usr/share/backuppc/html/0010000.gif
usr/share/backuppc/html/0010001.gif
usr/share/backuppc/html/0011000.gif
usr/share/backuppc/html/0011001.gif
usr/share/backuppc/html/1000000.gif
usr/share/backuppc/html/1000100.gif
usr/share/backuppc/html/1001000.gif
usr/share/backuppc/html/1001100.gif
usr/share/backuppc/html/1010000.gif
usr/share/backuppc/html/1010001.gif
usr/share/backuppc/html/1011000.gif
usr/share/backuppc/html/1100000.gif
usr/share/backuppc/html/1100100.gif
usr/share/backuppc/html/1100101.gif
usr/share/backuppc/html/1100110.gif
usr/share/backuppc/html/1100111.gif
usr/share/backuppc/html/1101000.gif
usr/share/backuppc/html/1101100.gif
usr/share/backuppc/html/1101101.gif
usr/share/backuppc/html/1101110.gif
usr/share/backuppc/html/1101111.gif
usr/share/backuppc/html/1110000.gif
usr/share/backuppc/html/1110001.gif
usr/share/backuppc/html/1110100.gif
usr/share/backuppc/html/1110101.gif
usr/share/backuppc/html/1110110.gif
usr/share/backuppc/html/1110111.gif
usr/share/backuppc/html/1111000.gif
usr/share/backuppc/html/1111001.gif
usr/share/backuppc/html/1111100.gif
usr/share/backuppc/html/1111101.gif
usr/share/backuppc/html/1111110.gif
usr/share/backuppc/html/1111111.gif
usr/share/backuppc/html/BackupPC_retro_v2.css
usr/share/backuppc/html/BackupPC_retro_v3.css
usr/share/backuppc/html/BackupPC_stnd.css
usr/share/backuppc/html/apple-touch-icon.png
usr/share/backuppc/html/favicon-16x16.png
usr/share/backuppc/html/favicon-32x32.png
usr/share/backuppc/html/favicon.ico
usr/share/backuppc/html/icon-dir.png
usr/share/backuppc/html/icon-file.png
usr/share/backuppc/html/icon-hardlink.png
usr/share/backuppc/html/icon-symlink.png
usr/share/backuppc/html/logo320.png
usr/share/backuppc/html/safari-pinned-tab.svg
usr/share/backuppc/html/sorttable.js
usr/share/backuppc/lib/
usr/share/backuppc/lib/BackupPC/
usr/share/backuppc/lib/BackupPC/CGI/
usr/share/backuppc/lib/BackupPC/CGI/AdminOptions.pm
usr/share/backuppc/lib/BackupPC/CGI/Archive.pm
usr/share/backuppc/lib/BackupPC/CGI/ArchiveInfo.pm
usr/share/backuppc/lib/BackupPC/CGI/Browse.pm
usr/share/backuppc/lib/BackupPC/CGI/DeleteBackup.pm
usr/share/backuppc/lib/BackupPC/CGI/DirHistory.pm
usr/share/backuppc/lib/BackupPC/CGI/EditConfig.pm
usr/share/backuppc/lib/BackupPC/CGI/EmailSummary.pm
usr/share/backuppc/lib/BackupPC/CGI/GeneralInfo.pm
usr/share/backuppc/lib/BackupPC/CGI/HostInfo.pm
usr/share/backuppc/lib/BackupPC/CGI/LOGlist.pm
usr/share/backuppc/lib/BackupPC/CGI/Lib.pm
usr/share/backuppc/lib/BackupPC/CGI/Metrics.pm
usr/share/backuppc/lib/BackupPC/CGI/Queue.pm
usr/share/backuppc/lib/BackupPC/CGI/ReloadServer.pm
usr/share/backuppc/lib/BackupPC/CGI/Restore.pm
usr/share/backuppc/lib/BackupPC/CGI/RestoreFile.pm
usr/share/backuppc/lib/BackupPC/CGI/RestoreInfo.pm
usr/share/backuppc/lib/BackupPC/CGI/StartServer.pm
usr/share/backuppc/lib/BackupPC/CGI/StartStopBackup.pm
usr/share/backuppc/lib/BackupPC/CGI/StopServer.pm
usr/share/backuppc/lib/BackupPC/CGI/Summary.pm
usr/share/backuppc/lib/BackupPC/CGI/View.pm
usr/share/backuppc/lib/BackupPC/Config/
usr/share/backuppc/lib/BackupPC/Config/Meta.pm
usr/share/backuppc/lib/BackupPC/DirOps.pm
usr/share/backuppc/lib/BackupPC/Lang/
usr/share/backuppc/lib/BackupPC/Lang/cz.pm
usr/share/backuppc/lib/BackupPC/Lang/de.pm
usr/share/backuppc/lib/BackupPC/Lang/en.pm
usr/share/backuppc/lib/BackupPC/Lang/es.pm
usr/share/backuppc/lib/BackupPC/Lang/fr.pm
usr/share/backuppc/lib/BackupPC/Lang/it.pm
usr/share/backuppc/lib/BackupPC/Lang/ja.pm
usr/share/backuppc/lib/BackupPC/Lang/nl.pm
usr/share/backuppc/lib/BackupPC/Lang/pl.pm
usr/share/backuppc/lib/BackupPC/Lang/pt_br.pm
usr/share/backuppc/lib/BackupPC/Lang/ru.pm
usr/share/backuppc/lib/BackupPC/Lang/uk.pm
usr/share/backuppc/lib/BackupPC/Lang/zh_CN.pm
usr/share/backuppc/lib/BackupPC/Lib.pm
usr/share/backuppc/lib/BackupPC/Storage.pm
usr/share/backuppc/lib/BackupPC/Storage/
usr/share/backuppc/lib/BackupPC/Storage/Text.pm
usr/share/backuppc/lib/BackupPC/View.pm
usr/share/backuppc/lib/BackupPC/Xfer.pm
usr/share/backuppc/lib/BackupPC/Xfer/
usr/share/backuppc/lib/BackupPC/Xfer/Archive.pm
usr/share/backuppc/lib/BackupPC/Xfer/Ftp.pm
usr/share/backuppc/lib/BackupPC/Xfer/Protocol.pm
usr/share/backuppc/lib/BackupPC/Xfer/Rsync.pm
usr/share/backuppc/lib/BackupPC/Xfer/Smb.pm
usr/share/backuppc/lib/BackupPC/Xfer/Tar.pm
usr/share/backuppc/lib/BackupPC/Zip/
usr/share/backuppc/lib/BackupPC/Zip/FileMember.pm
usr/share/backuppc/lib/Net/
usr/share/backuppc/lib/Net/FTP/
usr/share/backuppc/lib/Net/FTP/AutoReconnect.pm
usr/share/backuppc/lib/Net/FTP/RetrHandle.pm
usr/share/backuppc/share/
usr/share/backuppc/share/doc/
usr/share/backuppc/share/doc/BackupPC/
usr/share/backuppc/share/doc/BackupPC/BackupPC.html
usr/share/man/
usr/share/man/man3/
usr/share/man/man3/BackupPC::XS.3pm.gz
var/
var/lib/
var/lib/backuppc/
var/lib/backuppc/cpool/
var/lib/backuppc/pc/
var/lib/backuppc/pool/
var/log/
var/log/backuppc/
