# extracted/soapyremote-0.5.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/SoapySDRServer
usr/lib/
usr/lib/SoapySDR/
usr/lib/SoapySDR/modules0.8/
usr/lib/SoapySDR/modules0.8/libremoteSupport.so
usr/lib/sysctl.d/
usr/lib/sysctl.d/SoapySDRServer.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/SoapySDRServer.service
usr/share/
usr/share/licenses/
usr/share/licenses/soapyremote/
usr/share/licenses/soapyremote/LICENSE.txt
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/SoapySDRServer.1.gz
