# extracted/buildkit-0.18.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/buildctl
usr/bin/buildkitd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/buildkit.service
usr/lib/systemd/system/buildkit.socket
usr/lib/systemd/user/
usr/lib/systemd/user/buildkit-proxy.service
usr/lib/systemd/user/buildkit-proxy.socket
usr/lib/systemd/user/buildkit.service
usr/share/
usr/share/doc/
usr/share/doc/buildkit/
usr/share/doc/buildkit/README.md
usr/share/doc/buildkit/docs/
usr/share/doc/buildkit/docs/annotations.md
usr/share/doc/buildkit/docs/attestations/
usr/share/doc/buildkit/docs/attestations/README.md
usr/share/doc/buildkit/docs/attestations/attestation-storage.md
usr/share/doc/buildkit/docs/attestations/sbom-protocol.md
usr/share/doc/buildkit/docs/attestations/sbom.md
usr/share/doc/buildkit/docs/attestations/slsa-definitions.md
usr/share/doc/buildkit/docs/attestations/slsa-provenance.md
usr/share/doc/buildkit/docs/build-repro.md
usr/share/doc/buildkit/docs/buildinfo.md
usr/share/doc/buildkit/docs/buildkitd.toml.md
usr/share/doc/buildkit/docs/cni-networking.md
usr/share/doc/buildkit/docs/deprecated.md
usr/share/doc/buildkit/docs/dev/
usr/share/doc/buildkit/docs/dev/README.md
usr/share/doc/buildkit/docs/dev/dockerfile-llb.md
usr/share/doc/buildkit/docs/dev/merge-diff.md
usr/share/doc/buildkit/docs/dev/remote-debugging.md
usr/share/doc/buildkit/docs/dev/request-lifecycle.md
usr/share/doc/buildkit/docs/dev/solver.md
usr/share/doc/buildkit/docs/experimental-syntaxes.md
usr/share/doc/buildkit/docs/freebsd.md
usr/share/doc/buildkit/docs/generate.go
usr/share/doc/buildkit/docs/images-readme.md
usr/share/doc/buildkit/docs/multi-platform.md
usr/share/doc/buildkit/docs/nydus.md
usr/share/doc/buildkit/docs/reference/
usr/share/doc/buildkit/docs/reference/buildctl.md
usr/share/doc/buildkit/docs/rootless.md
usr/share/doc/buildkit/docs/stargz-estargz.md
usr/share/doc/buildkit/docs/windows.md
