# extracted/prometheus-unbound-exporter-0.4.6-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-unbound-exporter
usr/
usr/bin/
usr/bin/prometheus-unbound-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-unbound-exporter.service
usr/share/
usr/share/doc/
usr/share/doc/prometheus-unbound-exporter/
usr/share/doc/prometheus-unbound-exporter/README.md
usr/share/licenses/
usr/share/licenses/prometheus-unbound-exporter/
usr/share/licenses/prometheus-unbound-exporter/LICENSE
