# extracted/lprint-1.3.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/lprint
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lprint.service
usr/share/
usr/share/licenses/
usr/share/licenses/lprint/
usr/share/licenses/lprint/LICENSE
usr/share/licenses/lprint/NOTICE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/lprint-add.1.gz
usr/share/man/man1/lprint-cancel.1.gz
usr/share/man/man1/lprint-default.1.gz
usr/share/man/man1/lprint-delete.1.gz
usr/share/man/man1/lprint-devices.1.gz
usr/share/man/man1/lprint-drivers.1.gz
usr/share/man/man1/lprint-jobs.1.gz
usr/share/man/man1/lprint-modify.1.gz
usr/share/man/man1/lprint-options.1.gz
usr/share/man/man1/lprint-printers.1.gz
usr/share/man/man1/lprint-server.1.gz
usr/share/man/man1/lprint-shutdown.1.gz
usr/share/man/man1/lprint-status.1.gz
usr/share/man/man1/lprint-submit.1.gz
usr/share/man/man1/lprint.1.gz
usr/share/man/man5/
usr/share/man/man5/lprint.conf.5.gz
