# extracted/home-assistant-1:2024.12.3-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/home-assistant.service
usr/share/
usr/share/home-assistant/
usr/share/home-assistant/homeassistant-2024.12.3-py3-none-any.whl
