# extracted/smbnetfs-0.6.3-6/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/smbnetfs
etc/smbnetfs/
etc/smbnetfs/.smb/
etc/smbnetfs/.smb/smb.conf
etc/smbnetfs/.smb/smbnetfs.auth
etc/smbnetfs/.smb/smbnetfs.conf
etc/smbnetfs/.smb/smbnetfs.host
usr/
usr/bin/
usr/bin/smbnetfs
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/smbnetfs.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/smbnetfs.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/smbnetfs.conf
usr/share/
usr/share/doc/
usr/share/doc/smbnetfs/
usr/share/doc/smbnetfs/ENGLISH.FAQ
usr/share/doc/smbnetfs/INSTALL
usr/share/doc/smbnetfs/README
usr/share/doc/smbnetfs/RUSSIAN.FAQ
usr/share/doc/smbnetfs/smbnetfs.conf
