# extracted/qemu-vmsr-helper-9.1.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/qemu-vmsr-helper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/qemu-vmsr-helper.service
usr/lib/systemd/system/qemu-vmsr-helper.socket
usr/share/
usr/share/licenses/
usr/share/licenses/qemu-vmsr-helper/
usr/share/licenses/qemu-vmsr-helper/BSD-2-Clause.txt
usr/share/licenses/qemu-vmsr-helper/LICENSE
usr/share/licenses/qemu-vmsr-helper/audio.LICENSE.ISC.txt
usr/share/licenses/qemu-vmsr-helper/edk2-licenses.txt
usr/share/licenses/qemu-vmsr-helper/fuse.LICENSE.BSD-2-Clause.txt
usr/share/licenses/qemu-vmsr-helper/tcg-arm.LICENSE.MIT.txt
usr/share/licenses/qemu-vmsr-helper/tcg.LICENSE.MIT.txt
usr/share/licenses/qemu-vmsr-helper/tci.LICENSE.MIT.txt
