# extracted/nomad-1.9.3-1/files
%FILES%
etc/
etc/nomad.d/
etc/nomad.d/defaults.hcl
usr/
usr/bin/
usr/bin/nomad
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nomad.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nomad.conf
usr/share/
usr/share/licenses/
usr/share/licenses/nomad/
usr/share/licenses/nomad/LICENSE
