# extracted/influxdb-2.7.10-1/files
%FILES%
usr/
usr/bin/
usr/bin/influxd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/influxdb.service
usr/share/
usr/share/licenses/
usr/share/licenses/influxdb/
usr/share/licenses/influxdb/LICENSE
