# extracted/dkfilter-0.11-14/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/dkfilter
usr/
usr/bin/
usr/bin/dk_responder.pl
usr/bin/dkfilter.in
usr/bin/dkfilter.out
usr/bin/dkverify.pl
usr/lib/
usr/lib/DKMessage.pm
usr/lib/MSDW/
usr/lib/MSDW/SMTP/
usr/lib/MSDW/SMTP/Client.pm
usr/lib/MSDW/SMTP/Server.pm
usr/lib/Mail/
usr/lib/Mail/DomainKeys.pm
usr/lib/Mail/DomainKeys/
usr/lib/Mail/DomainKeys/Header.pm
usr/lib/Mail/DomainKeys/Key.pm
usr/lib/Mail/DomainKeys/Key/
usr/lib/Mail/DomainKeys/Key/Private.pm
usr/lib/Mail/DomainKeys/Key/Public.pm
usr/lib/Mail/DomainKeys/Message.pm
usr/lib/Mail/DomainKeys/Policy.pm
usr/lib/Mail/DomainKeys/Signature.pm
usr/lib/MySmtpProxyServer.pm
usr/lib/MySmtpServer.pm
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dkfilter-in.service
usr/lib/systemd/system/dkfilter-out.service
