# extracted/virtualbox-7.1.4-3/files
%FILES%
etc/
etc/vbox/
etc/vbox/vbox.cfg
usr/
usr/bin/
usr/bin/VBox
usr/bin/VBoxAutostart
usr/bin/VBoxBalloonCtrl
usr/bin/VBoxBugReport
usr/bin/VBoxHeadless
usr/bin/VBoxManage
usr/bin/VBoxSDL
usr/bin/VirtualBox
usr/bin/VirtualBoxVM
usr/bin/rcvboxdrv
usr/bin/vbox-img
usr/bin/vboxautostart
usr/bin/vboxballoonctrl
usr/bin/vboxbugreport
usr/bin/vboxheadless
usr/bin/vboximg-mount
usr/bin/vboxmanage
usr/bin/vboxreload
usr/bin/vboxsdl
usr/bin/vboxwebsrv
usr/bin/virtualbox
usr/bin/virtualboxvm
usr/lib/
usr/lib/modprobe.d/
usr/lib/modprobe.d/virtualbox.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vboxweb.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/virtualbox.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-vboxdrv.rules
usr/lib/virtualbox/
usr/lib/virtualbox/DbgPlugInDiggers.so
usr/lib/virtualbox/ExtensionPacks/
usr/lib/virtualbox/UICommon.so
usr/lib/virtualbox/VBoxAuth.so
usr/lib/virtualbox/VBoxAuthSimple.so
usr/lib/virtualbox/VBoxBalloonCtrl
usr/lib/virtualbox/VBoxDD.so
usr/lib/virtualbox/VBoxDD2.so
usr/lib/virtualbox/VBoxDDR0.r0
usr/lib/virtualbox/VBoxDDU.so
usr/lib/virtualbox/VBoxDbg.so
usr/lib/virtualbox/VBoxDragAndDropSvc.so
usr/lib/virtualbox/VBoxDxVk.so
usr/lib/virtualbox/VBoxEFI32.fd
usr/lib/virtualbox/VBoxEFI64.fd
usr/lib/virtualbox/VBoxExtPackHelperApp
usr/lib/virtualbox/VBoxGuestControlSvc.so
usr/lib/virtualbox/VBoxGuestPropSvc.so
usr/lib/virtualbox/VBoxHeadless
usr/lib/virtualbox/VBoxHeadless.so
usr/lib/virtualbox/VBoxHostChannel.so
usr/lib/virtualbox/VBoxKeyboard.so
usr/lib/virtualbox/VBoxManage
usr/lib/virtualbox/VBoxNetAdpCtl
usr/lib/virtualbox/VBoxNetDHCP
usr/lib/virtualbox/VBoxNetDHCP.so
usr/lib/virtualbox/VBoxNetNAT
usr/lib/virtualbox/VBoxNetNAT.so
usr/lib/virtualbox/VBoxRT.so
usr/lib/virtualbox/VBoxSDL
usr/lib/virtualbox/VBoxSDL.so
usr/lib/virtualbox/VBoxSVC
usr/lib/virtualbox/VBoxSVGA3D.so
usr/lib/virtualbox/VBoxSharedClipboard.so
usr/lib/virtualbox/VBoxSharedFolders.so
usr/lib/virtualbox/VBoxTraceLogDecoders.so
usr/lib/virtualbox/VBoxVMM.so
usr/lib/virtualbox/VBoxVMMPreload.so
usr/lib/virtualbox/VBoxXPCOM.so
usr/lib/virtualbox/VBoxXPCOMC.so
usr/lib/virtualbox/VBoxXPCOMIPCD.so
usr/lib/virtualbox/VMMR0.r0
usr/lib/virtualbox/VirtualBox
usr/lib/virtualbox/VirtualBoxVM
usr/lib/virtualbox/VirtualBoxVM.so
usr/lib/virtualbox/components/
usr/lib/virtualbox/components/VBoxC.so
usr/lib/virtualbox/components/VBoxDDU.so
usr/lib/virtualbox/components/VBoxRT.so
usr/lib/virtualbox/components/VBoxSVCM.so
usr/lib/virtualbox/components/VBoxVMM.so
usr/lib/virtualbox/components/VBoxXPCOM.so
usr/lib/virtualbox/components/VBoxXPCOMBase.xpt
usr/lib/virtualbox/components/VBoxXPCOMIPCC.so
usr/lib/virtualbox/components/VirtualBox_XPCOM.xpt
usr/lib/virtualbox/libvboxjxpcom.so
usr/lib/virtualbox/vbox-img
usr/lib/virtualbox/vboximg-mount
usr/lib/virtualbox/vboxwebsrv
usr/lib/virtualbox/webtest
usr/share/
usr/share/applications/
usr/share/applications/virtualbox.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/128x128/
usr/share/icons/hicolor/128x128/mimetypes/
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/128x128/mimetypes/virtualbox.png
usr/share/icons/hicolor/16x16/
usr/share/icons/hicolor/16x16/mimetypes/
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/16x16/mimetypes/virtualbox.png
usr/share/icons/hicolor/20x20/
usr/share/icons/hicolor/20x20/mimetypes/
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/20x20/mimetypes/virtualbox.png
usr/share/icons/hicolor/24x24/
usr/share/icons/hicolor/24x24/mimetypes/
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/24x24/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/256x256/
usr/share/icons/hicolor/256x256/mimetypes/
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/256x256/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/32x32/
usr/share/icons/hicolor/32x32/mimetypes/
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/32x32/mimetypes/virtualbox.png
usr/share/icons/hicolor/40x40/
usr/share/icons/hicolor/40x40/mimetypes/
usr/share/icons/hicolor/40x40/mimetypes/virtualbox.png
usr/share/icons/hicolor/48x48/
usr/share/icons/hicolor/48x48/mimetypes/
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/48x48/mimetypes/virtualbox.png
usr/share/icons/hicolor/64x64/
usr/share/icons/hicolor/64x64/mimetypes/
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/64x64/mimetypes/virtualbox.png
usr/share/icons/hicolor/72x72/
usr/share/icons/hicolor/72x72/mimetypes/
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/72x72/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/96x96/
usr/share/icons/hicolor/96x96/mimetypes/
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-hdd.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-ova.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-ovf.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-vbox-extpack.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-vbox.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-vdi.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-vhd.png
usr/share/icons/hicolor/96x96/mimetypes/virtualbox-vmdk.png
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/mimetypes/
usr/share/icons/hicolor/scalable/mimetypes/virtualbox.svg
usr/share/licenses/
usr/share/licenses/virtualbox/
usr/share/licenses/virtualbox/LICENSE
usr/share/licenses/virtualbox/LICENSE.CDDL
usr/share/mime/
usr/share/mime/packages/
usr/share/mime/packages/virtualbox.xml
usr/share/pixmaps/
usr/share/pixmaps/VBox.png
usr/share/virtualbox/
usr/share/virtualbox/UnattendedTemplates/
usr/share/virtualbox/UnattendedTemplates/debian_postinstall.sh
usr/share/virtualbox/UnattendedTemplates/debian_preseed.cfg
usr/share/virtualbox/UnattendedTemplates/fedora_ks.cfg
usr/share/virtualbox/UnattendedTemplates/freebsd_installer.cfg
usr/share/virtualbox/UnattendedTemplates/freebsd_postinstall.sh
usr/share/virtualbox/UnattendedTemplates/ol8_ks.cfg
usr/share/virtualbox/UnattendedTemplates/ol9_ks.cfg
usr/share/virtualbox/UnattendedTemplates/ol_ks.cfg
usr/share/virtualbox/UnattendedTemplates/ol_postinstall.sh
usr/share/virtualbox/UnattendedTemplates/os2_cid_install.cmd
usr/share/virtualbox/UnattendedTemplates/os2_response_files.rsp
usr/share/virtualbox/UnattendedTemplates/os2_util.exe
usr/share/virtualbox/UnattendedTemplates/redhat67_ks.cfg
usr/share/virtualbox/UnattendedTemplates/redhat_postinstall.sh
usr/share/virtualbox/UnattendedTemplates/rhel3_ks.cfg
usr/share/virtualbox/UnattendedTemplates/rhel4_ks.cfg
usr/share/virtualbox/UnattendedTemplates/rhel5_ks.cfg
usr/share/virtualbox/UnattendedTemplates/suse_autoinstall.xml
usr/share/virtualbox/UnattendedTemplates/ubuntu_autoinstall_meta_data
usr/share/virtualbox/UnattendedTemplates/ubuntu_autoinstall_user_data
usr/share/virtualbox/UnattendedTemplates/ubuntu_preseed.cfg
usr/share/virtualbox/UnattendedTemplates/win_nt5_unattended.sif
usr/share/virtualbox/UnattendedTemplates/win_nt6_unattended.xml
usr/share/virtualbox/UnattendedTemplates/win_postinstall.cmd
usr/share/virtualbox/VBoxCreateUSBNode.sh
usr/share/virtualbox/VBoxSysInfo.sh
usr/share/virtualbox/nls/
usr/share/virtualbox/nls/VBoxManageNls_ru.qm
usr/share/virtualbox/nls/VirtualBoxAPI_ru.qm
usr/share/virtualbox/nls/VirtualBoxAPI_tr.qm
usr/share/virtualbox/nls/VirtualBox_bg.qm
usr/share/virtualbox/nls/VirtualBox_ca.qm
usr/share/virtualbox/nls/VirtualBox_ca_VA.qm
usr/share/virtualbox/nls/VirtualBox_cs.qm
usr/share/virtualbox/nls/VirtualBox_da.qm
usr/share/virtualbox/nls/VirtualBox_de.qm
usr/share/virtualbox/nls/VirtualBox_el.qm
usr/share/virtualbox/nls/VirtualBox_en.qm
usr/share/virtualbox/nls/VirtualBox_es.qm
usr/share/virtualbox/nls/VirtualBox_eu.qm
usr/share/virtualbox/nls/VirtualBox_fa.qm
usr/share/virtualbox/nls/VirtualBox_fi.qm
usr/share/virtualbox/nls/VirtualBox_fr.qm
usr/share/virtualbox/nls/VirtualBox_gl.qm
usr/share/virtualbox/nls/VirtualBox_he.qm
usr/share/virtualbox/nls/VirtualBox_hr_HR.qm
usr/share/virtualbox/nls/VirtualBox_hu.qm
usr/share/virtualbox/nls/VirtualBox_id.qm
usr/share/virtualbox/nls/VirtualBox_it.qm
usr/share/virtualbox/nls/VirtualBox_ja.qm
usr/share/virtualbox/nls/VirtualBox_ka.qm
usr/share/virtualbox/nls/VirtualBox_km_KH.qm
usr/share/virtualbox/nls/VirtualBox_ko.qm
usr/share/virtualbox/nls/VirtualBox_lt.qm
usr/share/virtualbox/nls/VirtualBox_nl.qm
usr/share/virtualbox/nls/VirtualBox_pl.qm
usr/share/virtualbox/nls/VirtualBox_pt.qm
usr/share/virtualbox/nls/VirtualBox_pt_BR.qm
usr/share/virtualbox/nls/VirtualBox_ro.qm
usr/share/virtualbox/nls/VirtualBox_ru.qm
usr/share/virtualbox/nls/VirtualBox_sk.qm
usr/share/virtualbox/nls/VirtualBox_sl.qm
usr/share/virtualbox/nls/VirtualBox_sr.qm
usr/share/virtualbox/nls/VirtualBox_sv.qm
usr/share/virtualbox/nls/VirtualBox_th.qm
usr/share/virtualbox/nls/VirtualBox_tr.qm
usr/share/virtualbox/nls/VirtualBox_uk.qm
usr/share/virtualbox/nls/VirtualBox_zh_CN.qm
usr/share/virtualbox/nls/VirtualBox_zh_TW.qm
usr/share/virtualbox/nls/qt_bg.qm
usr/share/virtualbox/nls/qt_ca.qm
usr/share/virtualbox/nls/qt_ca_VA.qm
usr/share/virtualbox/nls/qt_cs.qm
usr/share/virtualbox/nls/qt_da.qm
usr/share/virtualbox/nls/qt_de.qm
usr/share/virtualbox/nls/qt_el.qm
usr/share/virtualbox/nls/qt_en.qm
usr/share/virtualbox/nls/qt_es.qm
usr/share/virtualbox/nls/qt_eu.qm
usr/share/virtualbox/nls/qt_fa.qm
usr/share/virtualbox/nls/qt_fi.qm
usr/share/virtualbox/nls/qt_fr.qm
usr/share/virtualbox/nls/qt_gl.qm
usr/share/virtualbox/nls/qt_he.qm
usr/share/virtualbox/nls/qt_hr_HR.qm
usr/share/virtualbox/nls/qt_hu.qm
usr/share/virtualbox/nls/qt_id.qm
usr/share/virtualbox/nls/qt_it.qm
usr/share/virtualbox/nls/qt_ja.qm
usr/share/virtualbox/nls/qt_ka.qm
usr/share/virtualbox/nls/qt_km_KH.qm
usr/share/virtualbox/nls/qt_ko.qm
usr/share/virtualbox/nls/qt_lt.qm
usr/share/virtualbox/nls/qt_nl.qm
usr/share/virtualbox/nls/qt_pl.qm
usr/share/virtualbox/nls/qt_pt.qm
usr/share/virtualbox/nls/qt_pt_BR.qm
usr/share/virtualbox/nls/qt_ro.qm
usr/share/virtualbox/nls/qt_ru.qm
usr/share/virtualbox/nls/qt_sk.qm
usr/share/virtualbox/nls/qt_sl.qm
usr/share/virtualbox/nls/qt_sr.qm
usr/share/virtualbox/nls/qt_sv.qm
usr/share/virtualbox/nls/qt_th.qm
usr/share/virtualbox/nls/qt_tr.qm
usr/share/virtualbox/nls/qt_uk.qm
usr/share/virtualbox/nls/qt_zh_CN.qm
usr/share/virtualbox/nls/qt_zh_TW.qm
