# extracted/charm-0.12.6-2/files
%FILES%
etc/
etc/charm.conf
usr/
usr/bin/
usr/bin/charm
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/charm.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/charm.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/charm.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/charm
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/charm.fish
usr/share/licenses/
usr/share/licenses/charm/
usr/share/licenses/charm/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/charm.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_charm
