# extracted/prometheus-node-exporter-1.8.2-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-node-exporter
usr/
usr/bin/
usr/bin/prometheus-node-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-node-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-node-exporter.conf
