# extracted/earlyoom-1.8.2-1/files
%FILES%
etc/
etc/default/
etc/default/earlyoom
usr/
usr/bin/
usr/bin/earlyoom
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/earlyoom.service
usr/share/
usr/share/licenses/
usr/share/licenses/earlyoom/
usr/share/licenses/earlyoom/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/earlyoom.1.gz
