# extracted/knockd-0.8-1/files
%FILES%
etc/
etc/knockd.conf
etc/logrotate.d/
etc/logrotate.d/knockd
usr/
usr/bin/
usr/bin/knock
usr/bin/knock_helper_ipt.sh
usr/bin/knockd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/knockd.service
usr/share/
usr/share/doc/
usr/share/doc/knock/
usr/share/doc/knock/COPYING
usr/share/doc/knock/ChangeLog
usr/share/doc/knock/README.md
usr/share/doc/knock/TODO
usr/share/doc/knockd/
usr/share/doc/knockd/README
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/knock.1.gz
usr/share/man/man1/knockd.1.gz
