# extracted/osquery-5.14.1-2/files
%FILES%
etc/
etc/osquery/
etc/osquery/osquery.conf
etc/sysconfig/
etc/sysconfig/osqueryd
usr/
usr/bin/
usr/bin/osqueryctl
usr/bin/osqueryd
usr/bin/osqueryi
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/osqueryd.service
usr/share/
usr/share/licenses/
usr/share/licenses/osquery/
usr/share/licenses/osquery/LICENSE
usr/share/osquery/
usr/share/osquery/certs/
usr/share/osquery/certs/certs.pem
usr/share/osquery/osquery.example.conf
usr/share/osquery/packs/
usr/share/osquery/packs/hardware-monitoring.conf
usr/share/osquery/packs/incident-response.conf
usr/share/osquery/packs/it-compliance.conf
usr/share/osquery/packs/osquery-monitoring.conf
usr/share/osquery/packs/ossec-rootkit.conf
usr/share/osquery/packs/osx-attacks.conf
usr/share/osquery/packs/unwanted-chrome-extensions.conf
usr/share/osquery/packs/vuln-management.conf
usr/share/osquery/packs/windows-attacks.conf
usr/share/osquery/packs/windows-hardening.conf
var/
var/log/
var/log/osquery/
