# extracted/drone-runner-exec-1.0.0.beta.10-2/files
%FILES%
usr/
usr/bin/
usr/bin/drone-runner-exec
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/drone-runner-exec.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/drone-runner-exec.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/drone-runner-exec.conf
usr/share/
usr/share/doc/
usr/share/doc/drone-runner-exec/
usr/share/doc/drone-runner-exec/CHANGELOG.md
usr/share/doc/drone-runner-exec/README.md
usr/share/licenses/
usr/share/licenses/drone-runner-exec/
usr/share/licenses/drone-runner-exec/LICENSE.md
usr/share/licenses/drone-runner-exec/Polyform-Free-Trial.md
usr/share/licenses/drone-runner-exec/Polyform-Noncommercial.md
usr/share/licenses/drone-runner-exec/Polyform-Small-Business.md
