# extracted/moosefs-4.56.6-1/files
%FILES%
etc/
etc/mfs/
etc/mfs/mfschunkserver.cfg.sample
etc/mfs/mfsexports.cfg.sample
etc/mfs/mfshdd.cfg.sample
etc/mfs/mfsipmap.cfg.sample
etc/mfs/mfsmaster.cfg.sample
etc/mfs/mfsmetalogger.cfg.sample
etc/mfs/mfsmount.cfg.sample
etc/mfs/mfstopology.cfg.sample
usr/
usr/bin/
usr/bin/mfsappendchunks
usr/bin/mfsarchive
usr/bin/mfsbdev
usr/bin/mfscgiserv
usr/bin/mfscheckfile
usr/bin/mfschkarchive
usr/bin/mfschunkdbdump
usr/bin/mfschunkserver
usr/bin/mfschunktool
usr/bin/mfscli
usr/bin/mfsclonesclass
usr/bin/mfsclrarchive
usr/bin/mfscopyeattr
usr/bin/mfscopygoal
usr/bin/mfscopyquota
usr/bin/mfscopysclass
usr/bin/mfscopytrashretention
usr/bin/mfscopytrashtime
usr/bin/mfscreatepattern
usr/bin/mfscreatesclass
usr/bin/mfscsstatsdump
usr/bin/mfsdeleattr
usr/bin/mfsdeletepattern
usr/bin/mfsdeletesclass
usr/bin/mfsdelquota
usr/bin/mfsdiagtools
usr/bin/mfsdirinfo
usr/bin/mfseattr
usr/bin/mfsfacl
usr/bin/mfsfileinfo
usr/bin/mfsfilepaths
usr/bin/mfsfilerepair
usr/bin/mfsgeteattr
usr/bin/mfsgetfacl
usr/bin/mfsgetgoal
usr/bin/mfsgetquota
usr/bin/mfsgetsclass
usr/bin/mfsgettrashretention
usr/bin/mfsgettrashtime
usr/bin/mfsimportsclass
usr/bin/mfslistpattern
usr/bin/mfslistsclass
usr/bin/mfsmakesnapshot
usr/bin/mfsmaster
usr/bin/mfsmetadirinfo
usr/bin/mfsmetadump
usr/bin/mfsmetalogger
usr/bin/mfsmetarestore
usr/bin/mfsmetasearch
usr/bin/mfsmodifysclass
usr/bin/mfsmount
usr/bin/mfsnetdump
usr/bin/mfspatadmin
usr/bin/mfsquota
usr/bin/mfsrenamesclass
usr/bin/mfsrmsnapshot
usr/bin/mfsscadmin
usr/bin/mfssclass
usr/bin/mfssetarchive
usr/bin/mfsseteattr
usr/bin/mfssetfacl
usr/bin/mfssetgoal
usr/bin/mfssetquota
usr/bin/mfssetsclass
usr/bin/mfssettrashretention
usr/bin/mfssettrashtime
usr/bin/mfssnapshots
usr/bin/mfsstatsdump
usr/bin/mfssupervisor
usr/bin/mfstrashretention
usr/bin/mfstrashtime
usr/bin/mfstrashtool
usr/bin/mfsxchgsclass
usr/bin/mount.moosefs
usr/include/
usr/include/mfsio.h
usr/lib/
usr/lib/libmfsio.so
usr/lib/libmfsio.so.1
usr/lib/libmfsio.so.1.0.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/moosefs-cgiserv.service
usr/lib/systemd/system/moosefs-chunkserver.service
usr/lib/systemd/system/moosefs-chunkserver@.service
usr/lib/systemd/system/moosefs-master.service
usr/lib/systemd/system/moosefs-master@.service
usr/lib/systemd/system/moosefs-metalogger.service
usr/lib/systemd/system/moosefs-metalogger@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/moosefs.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/moosefs.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mfsappendchunks.1.gz
usr/share/man/man1/mfsarchive.1.gz
usr/share/man/man1/mfscheckfile.1.gz
usr/share/man/man1/mfschkarchive.1.gz
usr/share/man/man1/mfscli.1.gz
usr/share/man/man1/mfsclonesclass.1.gz
usr/share/man/man1/mfsclrarchive.1.gz
usr/share/man/man1/mfscopyeattr.1.gz
usr/share/man/man1/mfscopygoal.1.gz
usr/share/man/man1/mfscopyquota.1.gz
usr/share/man/man1/mfscopysclass.1.gz
usr/share/man/man1/mfscopytrashretention.1.gz
usr/share/man/man1/mfscopytrashtime.1.gz
usr/share/man/man1/mfscreatepattern.1.gz
usr/share/man/man1/mfscreatesclass.1.gz
usr/share/man/man1/mfsdeleattr.1.gz
usr/share/man/man1/mfsdeletepattern.1.gz
usr/share/man/man1/mfsdeletesclass.1.gz
usr/share/man/man1/mfsdelquota.1.gz
usr/share/man/man1/mfsdiagtools.1.gz
usr/share/man/man1/mfsdirinfo.1.gz
usr/share/man/man1/mfseattr.1.gz
usr/share/man/man1/mfsfacl.1.gz
usr/share/man/man1/mfsfileinfo.1.gz
usr/share/man/man1/mfsfilepaths.1.gz
usr/share/man/man1/mfsfilerepair.1.gz
usr/share/man/man1/mfsgeteattr.1.gz
usr/share/man/man1/mfsgetfacl.1.gz
usr/share/man/man1/mfsgetgoal.1.gz
usr/share/man/man1/mfsgetquota.1.gz
usr/share/man/man1/mfsgetsclass.1.gz
usr/share/man/man1/mfsgettrashretention.1.gz
usr/share/man/man1/mfsgettrashtime.1.gz
usr/share/man/man1/mfsgoal.1.gz
usr/share/man/man1/mfsimportsclass.1.gz
usr/share/man/man1/mfslistpattern.1.gz
usr/share/man/man1/mfslistsclass.1.gz
usr/share/man/man1/mfsmakesnapshot.1.gz
usr/share/man/man1/mfsmodifysclass.1.gz
usr/share/man/man1/mfspatadmin.1.gz
usr/share/man/man1/mfsquota.1.gz
usr/share/man/man1/mfsrenamesclass.1.gz
usr/share/man/man1/mfsrmsnapshot.1.gz
usr/share/man/man1/mfsscadmin.1.gz
usr/share/man/man1/mfssclass.1.gz
usr/share/man/man1/mfssetarchive.1.gz
usr/share/man/man1/mfsseteattr.1.gz
usr/share/man/man1/mfssetfacl.1.gz
usr/share/man/man1/mfssetgoal.1.gz
usr/share/man/man1/mfssetquota.1.gz
usr/share/man/man1/mfssetsclass.1.gz
usr/share/man/man1/mfssettrashretention.1.gz
usr/share/man/man1/mfssettrashtime.1.gz
usr/share/man/man1/mfssnapshots.1.gz
usr/share/man/man1/mfstools.1.gz
usr/share/man/man1/mfstrashretention.1.gz
usr/share/man/man1/mfstrashtime.1.gz
usr/share/man/man1/mfsxchgsclass.1.gz
usr/share/man/man5/
usr/share/man/man5/mfsbdev.cfg.5.gz
usr/share/man/man5/mfschunkserver.cfg.5.gz
usr/share/man/man5/mfsexports.cfg.5.gz
usr/share/man/man5/mfshdd.cfg.5.gz
usr/share/man/man5/mfsipmap.cfg.5.gz
usr/share/man/man5/mfsmaster.cfg.5.gz
usr/share/man/man5/mfsmetalogger.cfg.5.gz
usr/share/man/man5/mfstopology.cfg.5.gz
usr/share/man/man7/
usr/share/man/man7/moosefs.7.gz
usr/share/man/man8/
usr/share/man/man8/mfsbdev.8.gz
usr/share/man/man8/mfscgiserv.8.gz
usr/share/man/man8/mfschunkdbdump.8.gz
usr/share/man/man8/mfschunkserver.8.gz
usr/share/man/man8/mfschunktool.8.gz
usr/share/man/man8/mfscsstatsdump.8.gz
usr/share/man/man8/mfsmaster.8.gz
usr/share/man/man8/mfsmetadirinfo.8.gz
usr/share/man/man8/mfsmetadump.8.gz
usr/share/man/man8/mfsmetalogger.8.gz
usr/share/man/man8/mfsmetarestore.8.gz
usr/share/man/man8/mfsmetasearch.8.gz
usr/share/man/man8/mfsmount.8.gz
usr/share/man/man8/mfsnetdump.8.gz
usr/share/man/man8/mfsstatsdump.8.gz
usr/share/man/man8/mfssupervisor.8.gz
usr/share/man/man8/mount.moosefs.8.gz
usr/share/mfscgi/
usr/share/mfscgi/acidchart.js
usr/share/mfscgi/acidtab.js
usr/share/mfscgi/cgiscripts.js
usr/share/mfscgi/chart.cgi
usr/share/mfscgi/chartdata.cgi
usr/share/mfscgi/err.gif
usr/share/mfscgi/favicon.ico
usr/share/mfscgi/index.html
usr/share/mfscgi/logo.svg
usr/share/mfscgi/logopro.svg
usr/share/mfscgi/mfs.cgi
usr/share/mfscgi/mfs.css
usr/share/mfscgi/mfsgraph.py
var/
var/lib/
var/lib/mfs/
var/lib/mfs/metadata.mfs.empty
