# extracted/kubelet-1.31.3-1/files
%FILES%
etc/
etc/kubernetes/
etc/kubernetes/kubelet.env
etc/kubernetes/manifests/
etc/sysctl.d/
etc/sysctl.d/50-kubelet.conf
usr/
usr/bin/
usr/bin/kubelet
usr/lib/
usr/lib/modules-load.d/
usr/lib/modules-load.d/kubelet.conf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/kubelet.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/kubelet.1.gz
