# extracted/xorg-xdm-1.1.16-1/files
%FILES%
etc/
etc/X11/
etc/X11/xdm/
etc/X11/xdm/GiveConsole
etc/X11/xdm/TakeConsole
etc/X11/xdm/Xaccess
etc/X11/xdm/Xreset
etc/X11/xdm/Xresources
etc/X11/xdm/Xservers
etc/X11/xdm/Xsession
etc/X11/xdm/Xsetup_0
etc/X11/xdm/Xstartup
etc/X11/xdm/Xwilling
etc/X11/xdm/xdm-config
etc/pam.d/
etc/pam.d/xdm
usr/
usr/bin/
usr/bin/xdm
usr/lib/
usr/lib/X11/
usr/lib/X11/xdm/
usr/lib/X11/xdm/chooser
usr/lib/X11/xdm/libXdmGreet.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/xdm.service
usr/share/
usr/share/X11/
usr/share/X11/app-defaults/
usr/share/X11/app-defaults/Chooser
usr/share/licenses/
usr/share/licenses/xorg-xdm/
usr/share/licenses/xorg-xdm/COPYING
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/xdm.8.gz
usr/share/xdm/
usr/share/xdm/pixmaps/
usr/share/xdm/pixmaps/xorg-bw.xpm
usr/share/xdm/pixmaps/xorg.xpm
var/
var/lib/
var/lib/xdm/
