# extracted/inadyn-2.12.0-2/files
%FILES%
etc/
etc/default/
etc/default/inadyn
etc/inadyn.conf
usr/
usr/bin/
usr/bin/inadyn
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/inadyn.service
usr/lib/systemd/system/inadyn.service.d/
usr/lib/systemd/system/inadyn.service.d/override.conf
usr/lib/sysusers.d/
usr/lib/sysusers.d/inadyn.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/inadyn.conf
usr/share/
usr/share/doc/
usr/share/doc/inadyn/
usr/share/doc/inadyn/COPYING
usr/share/doc/inadyn/ChangeLog.md
usr/share/doc/inadyn/README.md
usr/share/doc/inadyn/examples/
usr/share/doc/inadyn/examples/README.md
usr/share/doc/inadyn/examples/cloudflare-ipv4-ipv6.conf
usr/share/doc/inadyn/examples/cloudflare-ipv4-only.conf
usr/share/doc/inadyn/examples/cloudflare-ipv6-only.conf
usr/share/doc/inadyn/examples/custom.conf
usr/share/doc/inadyn/examples/duckdns.conf
usr/share/doc/inadyn/examples/dyndns.conf
usr/share/doc/inadyn/examples/freedns.conf
usr/share/doc/inadyn/examples/freemyip.conf
usr/share/doc/inadyn/examples/inadyn.conf
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/inadyn.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/inadyn.8.gz
