# extracted/headscale-0.23.0-1/files
%FILES%
etc/
etc/headscale/
etc/headscale/config.yaml
etc/headscale/derp.yaml
usr/
usr/bin/
usr/bin/headscale
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/headscale.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/headscale.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/headscale.conf
usr/share/
usr/share/headscale/
usr/share/headscale/config-example.yaml
usr/share/headscale/derp-example.yaml
usr/share/licenses/
usr/share/licenses/headscale/
usr/share/licenses/headscale/LICENSE
