# extracted/usbctl-1.2-3/files
%FILES%
usr/
usr/bin/
usr/bin/usbctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/deny-new-usb.service
usr/share/
usr/share/bash/
usr/share/bash/bash-completion/
usr/share/bash/bash-completion/completions/
usr/share/bash/bash-completion/completions/usbctl
usr/share/licenses/
usr/share/licenses/usbctl/
usr/share/licenses/usbctl/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/usbctl.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_usbctl
