# extracted/system76-scheduler-2.0.2-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/com.system76.Scheduler.conf
etc/system76-scheduler/
etc/system76-scheduler/config.kdl
etc/system76-scheduler/process-scheduler/
etc/system76-scheduler/process-scheduler/pop_os.kdl
usr/
usr/bin/
usr/bin/system76-scheduler
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/com.system76.Scheduler.service
