# extracted/v2ray-5.18.0-1/files
%FILES%
etc/
etc/v2ray/
etc/v2ray/config.json
etc/v2ray/vpoint_socks_vmess.json
etc/v2ray/vpoint_vmess_freedom.json
usr/
usr/bin/
usr/bin/v2ray
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/v2ray.service
usr/lib/systemd/system/v2ray@.service
usr/share/
usr/share/licenses/
usr/share/licenses/v2ray/
usr/share/licenses/v2ray/LICENSE
