# extracted/zerotier-one-1.14.2-2/files
%FILES%
usr/
usr/bin/
usr/bin/zerotier-cli
usr/bin/zerotier-idtool
usr/bin/zerotier-one
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/zerotier-one.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/zerotier-one.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/zerotier-cli.1.gz
usr/share/man/man1/zerotier-idtool.1.gz
usr/share/man/man8/
usr/share/man/man8/zerotier-one.8.gz
var/
var/lib/
var/lib/zerotier-one/
var/lib/zerotier-one/zerotier-cli
var/lib/zerotier-one/zerotier-idtool
var/lib/zerotier-one/zerotier-one
var/lib/zerotier-one/zerotier-one.te
