# extracted/trojan-1.16.0-12/files
%FILES%
etc/
etc/trojan/
etc/trojan/config.json
usr/
usr/bin/
usr/bin/trojan
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/trojan.service
usr/lib/systemd/system/trojan@.service
usr/share/
usr/share/doc/
usr/share/doc/trojan/
usr/share/doc/trojan/README.md
usr/share/doc/trojan/authenticator.md
usr/share/doc/trojan/build.md
usr/share/doc/trojan/config.md
usr/share/doc/trojan/examples/
usr/share/doc/trojan/examples/client.json-example
usr/share/doc/trojan/examples/forward.json-example
usr/share/doc/trojan/examples/nat.json-example
usr/share/doc/trojan/examples/server.json-example
usr/share/doc/trojan/overview.md
usr/share/doc/trojan/protocol.md
usr/share/doc/trojan/usage.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/trojan.1.gz
