# extracted/prometheus-memcached-exporter-0.15.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-memcached-exporter
usr/
usr/bin/
usr/bin/prometheus-memcached-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-memcached-exporter.service
