# extracted/stalwart-mail-0.9.2-1/files
%FILES%
etc/
etc/stalwart-mail/
etc/stalwart-mail/config.toml
usr/
usr/bin/
usr/bin/stalwart-mail
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/stalwart-mail.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/stalwart-mail.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/stalwart-mail.conf
usr/share/
usr/share/doc/
usr/share/doc/stalwart-mail/
usr/share/doc/stalwart-mail/CHANGELOG.md
usr/share/doc/stalwart-mail/README.md
usr/share/doc/stalwart-mail/UPGRADING.md
usr/share/stalwart-mail/
usr/share/stalwart-mail/build.py
usr/share/stalwart-mail/spamfilter.toml
usr/share/stalwart-mail/spamfilter/
usr/share/stalwart-mail/spamfilter/maps/
usr/share/stalwart-mail/spamfilter/maps/allow_dmarc.list
usr/share/stalwart-mail/spamfilter/maps/allow_domains.list
usr/share/stalwart-mail/spamfilter/maps/allow_spf_dkim.list
usr/share/stalwart-mail/spamfilter/maps/domains_disposable.list
usr/share/stalwart-mail/spamfilter/maps/domains_free.list
usr/share/stalwart-mail/spamfilter/maps/mime_types.map
usr/share/stalwart-mail/spamfilter/maps/scores.map
usr/share/stalwart-mail/spamfilter/maps/spam_config.map
usr/share/stalwart-mail/spamfilter/maps/url_redirectors.list
usr/share/stalwart-mail/spamfilter/scripts/
usr/share/stalwart-mail/spamfilter/scripts/bayes_classify.sieve
usr/share/stalwart-mail/spamfilter/scripts/bounce.sieve
usr/share/stalwart-mail/spamfilter/scripts/composites.sieve
usr/share/stalwart-mail/spamfilter/scripts/config.sieve
usr/share/stalwart-mail/spamfilter/scripts/date.sieve
usr/share/stalwart-mail/spamfilter/scripts/dmarc.sieve
usr/share/stalwart-mail/spamfilter/scripts/epilogue.sieve
usr/share/stalwart-mail/spamfilter/scripts/from.sieve
usr/share/stalwart-mail/spamfilter/scripts/greylist.sieve
usr/share/stalwart-mail/spamfilter/scripts/headers.sieve
usr/share/stalwart-mail/spamfilter/scripts/helo.sieve
usr/share/stalwart-mail/spamfilter/scripts/html.sieve
usr/share/stalwart-mail/spamfilter/scripts/ip.sieve
usr/share/stalwart-mail/spamfilter/scripts/messageid.sieve
usr/share/stalwart-mail/spamfilter/scripts/mime.sieve
usr/share/stalwart-mail/spamfilter/scripts/prelude.sieve
usr/share/stalwart-mail/spamfilter/scripts/pyzor.sieve
usr/share/stalwart-mail/spamfilter/scripts/rbl.sieve
usr/share/stalwart-mail/spamfilter/scripts/received.sieve
usr/share/stalwart-mail/spamfilter/scripts/recipient.sieve
usr/share/stalwart-mail/spamfilter/scripts/replies_in.sieve
usr/share/stalwart-mail/spamfilter/scripts/replies_out.sieve
usr/share/stalwart-mail/spamfilter/scripts/replyto.sieve
usr/share/stalwart-mail/spamfilter/scripts/reputation.sieve
usr/share/stalwart-mail/spamfilter/scripts/scores.sieve
usr/share/stalwart-mail/spamfilter/scripts/spamtrap.sieve
usr/share/stalwart-mail/spamfilter/scripts/subject.sieve
usr/share/stalwart-mail/spamfilter/scripts/train.sieve
usr/share/stalwart-mail/spamfilter/scripts/url.sieve
