# extracted/geoipupdate-7.1.0-1/files
%FILES%
etc/
etc/GeoIP.conf
usr/
usr/bin/
usr/bin/geoipupdate
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/geoipupdate.service
usr/lib/systemd/system/geoipupdate.timer
usr/share/
usr/share/doc/
usr/share/doc/geoipupdate/
usr/share/doc/geoipupdate/GeoIP.conf.md
usr/share/doc/geoipupdate/geoipupdate.md
usr/share/licenses/
usr/share/licenses/geoipupdate/
usr/share/licenses/geoipupdate/LICENSE-APACHE
usr/share/licenses/geoipupdate/LICENSE-MIT
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/geoipupdate.1.gz
usr/share/man/man5/
usr/share/man/man5/GeoIP.conf.5.gz
var/
var/lib/
var/lib/GeoIP/
