# extracted/radicale-3.3.1-1/files
%FILES%
etc/
etc/radicale/
etc/radicale/config
etc/radicale/rights
etc/radicale/users
usr/
usr/bin/
usr/bin/radicale
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/COPYING.md
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/METADATA
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/RECORD
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/WHEEL
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/entry_points.txt
usr/lib/python3.12/site-packages/Radicale-3.3.1.dist-info/top_level.txt
usr/lib/python3.12/site-packages/radicale/
usr/lib/python3.12/site-packages/radicale/__init__.py
usr/lib/python3.12/site-packages/radicale/__main__.py
usr/lib/python3.12/site-packages/radicale/__pycache__/
usr/lib/python3.12/site-packages/radicale/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/__main__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/__main__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/config.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/config.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/httputils.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/httputils.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/log.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/log.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/pathutils.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/pathutils.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/server.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/server.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/types.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/types.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/utils.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/utils.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/xmlutils.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/__pycache__/xmlutils.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/
usr/lib/python3.12/site-packages/radicale/app/__init__.py
usr/lib/python3.12/site-packages/radicale/app/__pycache__/
usr/lib/python3.12/site-packages/radicale/app/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/base.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/base.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/delete.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/delete.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/get.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/get.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/head.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/head.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/mkcalendar.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/mkcalendar.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/mkcol.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/mkcol.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/move.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/move.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/options.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/options.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/post.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/post.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/propfind.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/propfind.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/proppatch.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/proppatch.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/put.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/put.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/report.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/app/__pycache__/report.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/app/base.py
usr/lib/python3.12/site-packages/radicale/app/delete.py
usr/lib/python3.12/site-packages/radicale/app/get.py
usr/lib/python3.12/site-packages/radicale/app/head.py
usr/lib/python3.12/site-packages/radicale/app/mkcalendar.py
usr/lib/python3.12/site-packages/radicale/app/mkcol.py
usr/lib/python3.12/site-packages/radicale/app/move.py
usr/lib/python3.12/site-packages/radicale/app/options.py
usr/lib/python3.12/site-packages/radicale/app/post.py
usr/lib/python3.12/site-packages/radicale/app/propfind.py
usr/lib/python3.12/site-packages/radicale/app/proppatch.py
usr/lib/python3.12/site-packages/radicale/app/put.py
usr/lib/python3.12/site-packages/radicale/app/report.py
usr/lib/python3.12/site-packages/radicale/auth/
usr/lib/python3.12/site-packages/radicale/auth/__init__.py
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/denyall.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/denyall.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/dovecot.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/dovecot.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/htpasswd.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/htpasswd.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/http_x_remote_user.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/http_x_remote_user.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/ldap.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/ldap.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/none.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/none.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/remote_user.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/auth/__pycache__/remote_user.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/auth/denyall.py
usr/lib/python3.12/site-packages/radicale/auth/dovecot.py
usr/lib/python3.12/site-packages/radicale/auth/htpasswd.py
usr/lib/python3.12/site-packages/radicale/auth/http_x_remote_user.py
usr/lib/python3.12/site-packages/radicale/auth/ldap.py
usr/lib/python3.12/site-packages/radicale/auth/none.py
usr/lib/python3.12/site-packages/radicale/auth/remote_user.py
usr/lib/python3.12/site-packages/radicale/config.py
usr/lib/python3.12/site-packages/radicale/hook/
usr/lib/python3.12/site-packages/radicale/hook/__init__.py
usr/lib/python3.12/site-packages/radicale/hook/__pycache__/
usr/lib/python3.12/site-packages/radicale/hook/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/hook/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/hook/__pycache__/none.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/hook/__pycache__/none.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/hook/none.py
usr/lib/python3.12/site-packages/radicale/hook/rabbitmq/
usr/lib/python3.12/site-packages/radicale/hook/rabbitmq/__init__.py
usr/lib/python3.12/site-packages/radicale/hook/rabbitmq/__pycache__/
usr/lib/python3.12/site-packages/radicale/hook/rabbitmq/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/hook/rabbitmq/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/httputils.py
usr/lib/python3.12/site-packages/radicale/item/
usr/lib/python3.12/site-packages/radicale/item/__init__.py
usr/lib/python3.12/site-packages/radicale/item/__pycache__/
usr/lib/python3.12/site-packages/radicale/item/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/item/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/item/__pycache__/filter.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/item/__pycache__/filter.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/item/filter.py
usr/lib/python3.12/site-packages/radicale/log.py
usr/lib/python3.12/site-packages/radicale/pathutils.py
usr/lib/python3.12/site-packages/radicale/py.typed
usr/lib/python3.12/site-packages/radicale/rights/
usr/lib/python3.12/site-packages/radicale/rights/__init__.py
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/authenticated.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/authenticated.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/from_file.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/from_file.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/owner_only.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/owner_only.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/owner_write.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/rights/__pycache__/owner_write.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/rights/authenticated.py
usr/lib/python3.12/site-packages/radicale/rights/from_file.py
usr/lib/python3.12/site-packages/radicale/rights/owner_only.py
usr/lib/python3.12/site-packages/radicale/rights/owner_write.py
usr/lib/python3.12/site-packages/radicale/server.py
usr/lib/python3.12/site-packages/radicale/storage/
usr/lib/python3.12/site-packages/radicale/storage/__init__.py
usr/lib/python3.12/site-packages/radicale/storage/__pycache__/
usr/lib/python3.12/site-packages/radicale/storage/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/__pycache__/multifilesystem_nolock.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/__pycache__/multifilesystem_nolock.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__init__.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/base.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/base.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/cache.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/cache.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/create_collection.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/create_collection.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/delete.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/delete.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/discover.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/discover.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/get.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/get.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/history.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/history.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/lock.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/lock.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/meta.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/meta.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/move.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/move.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/sync.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/sync.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/upload.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/upload.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/verify.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/__pycache__/verify.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/base.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/cache.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/create_collection.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/delete.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/discover.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/get.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/history.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/lock.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/meta.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/move.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/sync.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/upload.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem/verify.py
usr/lib/python3.12/site-packages/radicale/storage/multifilesystem_nolock.py
usr/lib/python3.12/site-packages/radicale/tests/
usr/lib/python3.12/site-packages/radicale/tests/custom/
usr/lib/python3.12/site-packages/radicale/tests/custom/__init__.py
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/auth.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/auth.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/rights.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/rights.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/storage_simple_sync.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/storage_simple_sync.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/web.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/__pycache__/web.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/tests/custom/auth.py
usr/lib/python3.12/site-packages/radicale/tests/custom/rights.py
usr/lib/python3.12/site-packages/radicale/tests/custom/storage_simple_sync.py
usr/lib/python3.12/site-packages/radicale/tests/custom/web.py
usr/lib/python3.12/site-packages/radicale/types.py
usr/lib/python3.12/site-packages/radicale/utils.py
usr/lib/python3.12/site-packages/radicale/web/
usr/lib/python3.12/site-packages/radicale/web/__init__.py
usr/lib/python3.12/site-packages/radicale/web/__pycache__/
usr/lib/python3.12/site-packages/radicale/web/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/web/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/web/__pycache__/internal.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/web/__pycache__/internal.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/web/__pycache__/none.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/radicale/web/__pycache__/none.cpython-312.pyc
usr/lib/python3.12/site-packages/radicale/web/internal.py
usr/lib/python3.12/site-packages/radicale/web/internal_data/
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icon.png
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/delete.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/download.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/edit.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/new.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/icons/upload.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/loading.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/logo.svg
usr/lib/python3.12/site-packages/radicale/web/internal_data/css/main.css
usr/lib/python3.12/site-packages/radicale/web/internal_data/fn.js
usr/lib/python3.12/site-packages/radicale/web/internal_data/index.html
usr/lib/python3.12/site-packages/radicale/web/none.py
usr/lib/python3.12/site-packages/radicale/xmlutils.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/radicale.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/radicale.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/radicale.conf
usr/share/
usr/share/doc/
usr/share/doc/radicale/
usr/share/doc/radicale/CHANGELOG.md
usr/share/doc/radicale/DOCUMENTATION.md
usr/share/doc/radicale/README.md
usr/share/radicale/
usr/share/radicale/radicale.wsgi
