# extracted/gitlab-gitaly-17.6.2-1/files
%FILES%
etc/
etc/gitlab-gitaly/
etc/gitlab-gitaly/config.toml
usr/
usr/bin/
usr/bin/gitaly
usr/bin/gitaly-backup
usr/bin/gitaly-blackbox
usr/bin/gitaly-debug
usr/bin/gitaly-git-http-backend-v2.46
usr/bin/gitaly-git-http-backend-v2.47
usr/bin/gitaly-git-remote-http-v2.46
usr/bin/gitaly-git-remote-http-v2.47
usr/bin/gitaly-git-v2.46
usr/bin/gitaly-git-v2.47
usr/bin/gitaly-wrapper
usr/bin/praefect
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gitlab-gitaly.service
usr/share/
usr/share/licenses/
usr/share/licenses/gitlab-gitaly/
usr/share/licenses/gitlab-gitaly/LICENSE
usr/share/webapps/
usr/share/webapps/gitlab-gitaly/
