# extracted/dhcp-helper-1.2-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/dhcp-helper
usr/
usr/bin/
usr/bin/dhcp-helper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dhcp-helper.service
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/dhcp-helper.8.gz
