# extracted/lxd-6.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/fuidshift
usr/bin/lxc
usr/bin/lxc-to-lxd
usr/bin/lxd
usr/bin/lxd-agent
usr/bin/lxd-benchmark
usr/bin/lxd-migrate
usr/bin/lxd-user
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lxd.service
usr/lib/systemd/system/lxd.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/lxd.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/lxc
usr/share/doc/
usr/share/doc/lxd/
usr/share/doc/lxd/Makefile
usr/share/doc/lxd/Makefile.sp
usr/share/doc/lxd/README.md
usr/share/doc/lxd/api-extensions.md
usr/share/doc/lxd/api.md
usr/share/doc/lxd/architectures.md
usr/share/doc/lxd/authentication.md
usr/share/doc/lxd/backup.md
usr/share/doc/lxd/cloud-init.md
usr/share/doc/lxd/clustering.md
usr/share/doc/lxd/conf.py
usr/share/doc/lxd/config-options.md
usr/share/doc/lxd/config_options_cheat_sheet.md
usr/share/doc/lxd/container-environment.md
usr/share/doc/lxd/contributing.md
usr/share/doc/lxd/custom_conf.py
usr/share/doc/lxd/daemon-behavior.md
usr/share/doc/lxd/database.md
usr/share/doc/lxd/debugging.md
usr/share/doc/lxd/dev-lxd.md
usr/share/doc/lxd/doc-cheat-sheet-myst.md
usr/share/doc/lxd/environment.md
usr/share/doc/lxd/events.md
usr/share/doc/lxd/explanation/
usr/share/doc/lxd/explanation/authorization.md
usr/share/doc/lxd/explanation/clustering.md
usr/share/doc/lxd/explanation/index.md
usr/share/doc/lxd/explanation/instance_config.md
usr/share/doc/lxd/explanation/instances.md
usr/share/doc/lxd/explanation/lxc_show_info.md
usr/share/doc/lxd/explanation/lxd_lxc.md
usr/share/doc/lxd/explanation/networks.md
usr/share/doc/lxd/explanation/performance_tuning.md
usr/share/doc/lxd/explanation/projects.md
usr/share/doc/lxd/explanation/security.md
usr/share/doc/lxd/explanation/storage.md
usr/share/doc/lxd/external_resources.md
usr/share/doc/lxd/faq.md
usr/share/doc/lxd/getting_started.md
usr/share/doc/lxd/guest-os-compatibility.md
usr/share/doc/lxd/howto/
usr/share/doc/lxd/howto/access_documentation.md
usr/share/doc/lxd/howto/access_ui.md
usr/share/doc/lxd/howto/benchmark_performance.md
usr/share/doc/lxd/howto/cluster_config_networks.md
usr/share/doc/lxd/howto/cluster_config_storage.md
usr/share/doc/lxd/howto/cluster_form.md
usr/share/doc/lxd/howto/cluster_groups.md
usr/share/doc/lxd/howto/cluster_manage.md
usr/share/doc/lxd/howto/cluster_manage_instance.md
usr/share/doc/lxd/howto/cluster_recover.md
usr/share/doc/lxd/howto/container_gpu_passthrough_with_docker.md
usr/share/doc/lxd/howto/disaster_recovery.md
usr/share/doc/lxd/howto/dqlite_troubleshoot.md
usr/share/doc/lxd/howto/grafana.md
usr/share/doc/lxd/howto/images_copy.md
usr/share/doc/lxd/howto/images_create.md
usr/share/doc/lxd/howto/images_manage.md
usr/share/doc/lxd/howto/images_profiles.md
usr/share/doc/lxd/howto/images_remote.md
usr/share/doc/lxd/howto/import_machines_to_instances.md
usr/share/doc/lxd/howto/index.md
usr/share/doc/lxd/howto/initialize.md
usr/share/doc/lxd/howto/instances_access_files.md
usr/share/doc/lxd/howto/instances_backup.md
usr/share/doc/lxd/howto/instances_configure.md
usr/share/doc/lxd/howto/instances_console.md
usr/share/doc/lxd/howto/instances_create.md
usr/share/doc/lxd/howto/instances_manage.md
usr/share/doc/lxd/howto/instances_routed_nic_vm.md
usr/share/doc/lxd/howto/instances_troubleshoot.md
usr/share/doc/lxd/howto/logs_loki.md
usr/share/doc/lxd/howto/lxc_alias.md
usr/share/doc/lxd/howto/migrate_from_lxc.md
usr/share/doc/lxd/howto/move_instances.md
usr/share/doc/lxd/howto/network_acls.md
usr/share/doc/lxd/howto/network_bgp.md
usr/share/doc/lxd/howto/network_bridge_firewalld.md
usr/share/doc/lxd/howto/network_bridge_resolved.md
usr/share/doc/lxd/howto/network_configure.md
usr/share/doc/lxd/howto/network_create.md
usr/share/doc/lxd/howto/network_forwards.md
usr/share/doc/lxd/howto/network_increase_bandwidth.md
usr/share/doc/lxd/howto/network_ipam.md
usr/share/doc/lxd/howto/network_load_balancers.md
usr/share/doc/lxd/howto/network_ovn_peers.md
usr/share/doc/lxd/howto/network_ovn_setup.md
usr/share/doc/lxd/howto/network_zones.md
usr/share/doc/lxd/howto/projects_confine.md
usr/share/doc/lxd/howto/projects_create.md
usr/share/doc/lxd/howto/projects_work.md
usr/share/doc/lxd/howto/server_configure.md
usr/share/doc/lxd/howto/server_expose.md
usr/share/doc/lxd/howto/snap.md
usr/share/doc/lxd/howto/storage_backup_volume.md
usr/share/doc/lxd/howto/storage_buckets.md
usr/share/doc/lxd/howto/storage_create_instance.md
usr/share/doc/lxd/howto/storage_move_volume.md
usr/share/doc/lxd/howto/storage_pools.md
usr/share/doc/lxd/howto/storage_volumes.md
usr/share/doc/lxd/howto/troubleshoot.md
usr/share/doc/lxd/image-handling.md
usr/share/doc/lxd/images.md
usr/share/doc/lxd/images/
usr/share/doc/lxd/images/UI/
usr/share/doc/lxd/images/UI/create_instance_ex1.png
usr/share/doc/lxd/images/UI/create_instance_ex2-2.png
usr/share/doc/lxd/images/UI/create_instance_ex2.png
usr/share/doc/lxd/images/UI/create_instance_ex3.png
usr/share/doc/lxd/images/UI/create_instance_ex4.png
usr/share/doc/lxd/images/UI/create_project.png
usr/share/doc/lxd/images/UI/create_restr_project1.png
usr/share/doc/lxd/images/UI/create_restr_project2.png
usr/share/doc/lxd/images/UI/limits_memory_example.png
usr/share/doc/lxd/images/UI/routed_nic_create_instance.png
usr/share/doc/lxd/images/UI/snapshots_cron.png
usr/share/doc/lxd/images/UI/snapshots_daily.png
usr/share/doc/lxd/images/application-vs-system-containers.svg
usr/share/doc/lxd/images/grafana_add_datasource.png
usr/share/doc/lxd/images/grafana_configure_datasource.png
usr/share/doc/lxd/images/grafana_configure_loki.png
usr/share/doc/lxd/images/grafana_configure_prometheus.png
usr/share/doc/lxd/images/grafana_dashboard_import.png
usr/share/doc/lxd/images/grafana_dashboard_select_datasource.png
usr/share/doc/lxd/images/grafana_instances.png
usr/share/doc/lxd/images/grafana_resources.png
usr/share/doc/lxd/images/grafana_select_loki.png
usr/share/doc/lxd/images/grafana_select_prometheus.png
usr/share/doc/lxd/images/grafana_welcome.png
usr/share/doc/lxd/images/ovn_networking_1.svg
usr/share/doc/lxd/images/ovn_networking_2.svg
usr/share/doc/lxd/images/tutorial/
usr/share/doc/lxd/images/tutorial/add_disk_device.png
usr/share/doc/lxd/images/tutorial/broken_terminal.png
usr/share/doc/lxd/images/tutorial/create_desktop_vm.png
usr/share/doc/lxd/images/tutorial/create_instance.png
usr/share/doc/lxd/images/tutorial/create_project.png
usr/share/doc/lxd/images/tutorial/create_vm.png
usr/share/doc/lxd/images/tutorial/desktop_console.png
usr/share/doc/lxd/images/tutorial/hello_world_desktop.png
usr/share/doc/lxd/images/tutorial/instance_summary.png
usr/share/doc/lxd/images/tutorial/instances.png
usr/share/doc/lxd/images/tutorial/resource_limits.png
usr/share/doc/lxd/images/tutorial/root_disk_size.png
usr/share/doc/lxd/images/tutorial/yaml_configuration.png
usr/share/doc/lxd/images/ui_console.png
usr/share/doc/lxd/images/ui_security_warning.png
usr/share/doc/lxd/images/ui_set_up_certificates.png
usr/share/doc/lxd/images/ui_set_up_existing_cert.png
usr/share/doc/lxd/images/virtual-machines-vs-system-containers.svg
usr/share/doc/lxd/index.md
usr/share/doc/lxd/installing.md
usr/share/doc/lxd/instance-exec.md
usr/share/doc/lxd/instances.md
usr/share/doc/lxd/internals.md
usr/share/doc/lxd/metadata.txt
usr/share/doc/lxd/metrics.md
usr/share/doc/lxd/migration.md
usr/share/doc/lxd/myconf.py
usr/share/doc/lxd/networks.md
usr/share/doc/lxd/operation.md
usr/share/doc/lxd/production-setup.md
usr/share/doc/lxd/profiles.md
usr/share/doc/lxd/projects.md
usr/share/doc/lxd/reference/
usr/share/doc/lxd/reference/cluster_member_config.md
usr/share/doc/lxd/reference/devices.md
usr/share/doc/lxd/reference/devices_disk.md
usr/share/doc/lxd/reference/devices_gpu.md
usr/share/doc/lxd/reference/devices_infiniband.md
usr/share/doc/lxd/reference/devices_nic.md
usr/share/doc/lxd/reference/devices_none.md
usr/share/doc/lxd/reference/devices_pci.md
usr/share/doc/lxd/reference/devices_proxy.md
usr/share/doc/lxd/reference/devices_tpm.md
usr/share/doc/lxd/reference/devices_unix_block.md
usr/share/doc/lxd/reference/devices_unix_char.md
usr/share/doc/lxd/reference/devices_unix_hotplug.md
usr/share/doc/lxd/reference/devices_usb.md
usr/share/doc/lxd/reference/image_format.md
usr/share/doc/lxd/reference/index.md
usr/share/doc/lxd/reference/instance_options.md
usr/share/doc/lxd/reference/instance_properties.md
usr/share/doc/lxd/reference/instance_units.md
usr/share/doc/lxd/reference/manpages.md
usr/share/doc/lxd/reference/manpages/
usr/share/doc/lxd/reference/manpages/readme.txt
usr/share/doc/lxd/reference/network_bridge.md
usr/share/doc/lxd/reference/network_macvlan.md
usr/share/doc/lxd/reference/network_ovn.md
usr/share/doc/lxd/reference/network_physical.md
usr/share/doc/lxd/reference/network_sriov.md
usr/share/doc/lxd/reference/networks.md
usr/share/doc/lxd/reference/permissions.md
usr/share/doc/lxd/reference/preseed_yaml_fields.md
usr/share/doc/lxd/reference/projects.md
usr/share/doc/lxd/reference/provided_metrics.md
usr/share/doc/lxd/reference/remote_image_servers.md
usr/share/doc/lxd/reference/server_settings.md
usr/share/doc/lxd/reference/standard_devices.md
usr/share/doc/lxd/reference/storage_btrfs.md
usr/share/doc/lxd/reference/storage_ceph.md
usr/share/doc/lxd/reference/storage_cephfs.md
usr/share/doc/lxd/reference/storage_cephobject.md
usr/share/doc/lxd/reference/storage_dir.md
usr/share/doc/lxd/reference/storage_drivers.md
usr/share/doc/lxd/reference/storage_lvm.md
usr/share/doc/lxd/reference/storage_powerflex.md
usr/share/doc/lxd/reference/storage_zfs.md
usr/share/doc/lxd/reference/uefi_variables.md
usr/share/doc/lxd/related_topics.yaml
usr/share/doc/lxd/remotes.md
usr/share/doc/lxd/requirements.md
usr/share/doc/lxd/rest-api.md
usr/share/doc/lxd/rest-api.yaml
usr/share/doc/lxd/restapi_landing.md
usr/share/doc/lxd/security.md
usr/share/doc/lxd/server.md
usr/share/doc/lxd/storage.md
usr/share/doc/lxd/substitutions.yaml
usr/share/doc/lxd/support.md
usr/share/doc/lxd/syscall-interception.md
usr/share/doc/lxd/tutorial/
usr/share/doc/lxd/tutorial/first_steps.md
usr/share/doc/lxd/tutorial/index.md
usr/share/doc/lxd/tutorial/ui.md
usr/share/doc/lxd/userns-idmap.md
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/lxc.fish
usr/share/licenses/
usr/share/licenses/lxd/
usr/share/licenses/lxd/LICENCE
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_lxc
var/
var/log/
var/log/lxd/
