# extracted/tinyssh-20240101-1/files
%FILES%
etc/
etc/tinyssh/
usr/
usr/bin/
usr/bin/tinyssh-convert
usr/bin/tinysshd
usr/bin/tinysshd-makekey
usr/bin/tinysshd-printkey
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tinyssh@.service
usr/lib/systemd/system/tinyssh@.socket
usr/lib/systemd/system/tinysshgenkeys.service
usr/share/
usr/share/licenses/
usr/share/licenses/tinyssh/
usr/share/licenses/tinyssh/LICENCE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/tinysshd-makekey.8.gz
usr/share/man/man8/tinysshd-printkey.8.gz
usr/share/man/man8/tinysshd.8.gz
usr/share/man/man8/tinysshnoneauthd.8.gz
