# extracted/mkosi-24.3-2/files
%FILES%
etc/
etc/mkosi-initrd/
usr/
usr/bin/
usr/bin/mkosi
usr/bin/mkosi-initrd
usr/lib/
usr/lib/kernel/
usr/lib/kernel/install.d/
usr/lib/kernel/install.d/50-mkosi.install
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/LICENSE
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/LICENSE.GPL2
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/METADATA
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/RECORD
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/WHEEL
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/entry_points.txt
usr/lib/python3.13/site-packages/mkosi-24.3.dist-info/top_level.txt
usr/lib/python3.13/site-packages/mkosi/
usr/lib/python3.13/site-packages/mkosi/__init__.py
usr/lib/python3.13/site-packages/mkosi/__main__.py
usr/lib/python3.13/site-packages/mkosi/__pycache__/
usr/lib/python3.13/site-packages/mkosi/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/__main__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/__main__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/archive.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/archive.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/burn.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/burn.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/completion.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/completion.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/config.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/config.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/context.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/context.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/kmod.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/kmod.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/log.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/log.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/manifest.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/manifest.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/mounts.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/mounts.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/pager.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/pager.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/partition.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/partition.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/qemu.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/qemu.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/run.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/run.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/sandbox.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/sandbox.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/tree.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/tree.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/types.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/types.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/user.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/user.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/util.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/util.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/versioncomp.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/versioncomp.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/vmspawn.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/__pycache__/vmspawn.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/archive.py
usr/lib/python3.13/site-packages/mkosi/burn.py
usr/lib/python3.13/site-packages/mkosi/completion.py
usr/lib/python3.13/site-packages/mkosi/config.py
usr/lib/python3.13/site-packages/mkosi/context.py
usr/lib/python3.13/site-packages/mkosi/distributions/
usr/lib/python3.13/site-packages/mkosi/distributions/__init__.py
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/alma.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/alma.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/arch.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/arch.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/centos.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/centos.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/custom.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/custom.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/debian.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/debian.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/fedora.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/fedora.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/mageia.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/mageia.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/openmandriva.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/openmandriva.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/opensuse.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/opensuse.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rhel.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rhel.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rhel_ubi.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rhel_ubi.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rocky.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/rocky.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/ubuntu.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/__pycache__/ubuntu.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/distributions/alma.py
usr/lib/python3.13/site-packages/mkosi/distributions/arch.py
usr/lib/python3.13/site-packages/mkosi/distributions/centos.py
usr/lib/python3.13/site-packages/mkosi/distributions/custom.py
usr/lib/python3.13/site-packages/mkosi/distributions/debian.py
usr/lib/python3.13/site-packages/mkosi/distributions/fedora.py
usr/lib/python3.13/site-packages/mkosi/distributions/mageia.py
usr/lib/python3.13/site-packages/mkosi/distributions/openmandriva.py
usr/lib/python3.13/site-packages/mkosi/distributions/opensuse.py
usr/lib/python3.13/site-packages/mkosi/distributions/rhel.py
usr/lib/python3.13/site-packages/mkosi/distributions/rhel_ubi.py
usr/lib/python3.13/site-packages/mkosi/distributions/rocky.py
usr/lib/python3.13/site-packages/mkosi/distributions/ubuntu.py
usr/lib/python3.13/site-packages/mkosi/initrd/
usr/lib/python3.13/site-packages/mkosi/initrd/__main__.py
usr/lib/python3.13/site-packages/mkosi/initrd/__pycache__/
usr/lib/python3.13/site-packages/mkosi/initrd/__pycache__/__main__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/initrd/__pycache__/__main__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/initrd/resources/
usr/lib/python3.13/site-packages/mkosi/initrd/resources/__init__.py
usr/lib/python3.13/site-packages/mkosi/initrd/resources/__pycache__/
usr/lib/python3.13/site-packages/mkosi/initrd/resources/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/initrd/resources/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/initrd/resources/mkosi-initrd.md
usr/lib/python3.13/site-packages/mkosi/installer/
usr/lib/python3.13/site-packages/mkosi/installer/__init__.py
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/apt.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/apt.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/dnf.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/dnf.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/pacman.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/pacman.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/rpm.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/rpm.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/zypper.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/installer/__pycache__/zypper.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/installer/apt.py
usr/lib/python3.13/site-packages/mkosi/installer/dnf.py
usr/lib/python3.13/site-packages/mkosi/installer/pacman.py
usr/lib/python3.13/site-packages/mkosi/installer/rpm.py
usr/lib/python3.13/site-packages/mkosi/installer/zypper.py
usr/lib/python3.13/site-packages/mkosi/kmod.py
usr/lib/python3.13/site-packages/mkosi/log.py
usr/lib/python3.13/site-packages/mkosi/manifest.py
usr/lib/python3.13/site-packages/mkosi/mounts.py
usr/lib/python3.13/site-packages/mkosi/pager.py
usr/lib/python3.13/site-packages/mkosi/partition.py
usr/lib/python3.13/site-packages/mkosi/qemu.py
usr/lib/python3.13/site-packages/mkosi/resources/
usr/lib/python3.13/site-packages/mkosi/resources/__init__.py
usr/lib/python3.13/site-packages/mkosi/resources/__pycache__/
usr/lib/python3.13/site-packages/mkosi/resources/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/mkosi/resources/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-arch.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-centos-fedora.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-centos.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-debian-ubuntu/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-debian-ubuntu/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/10-libtss.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-fedora.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/10-opensuse.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.conf.d/20-stub.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/system-preset/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/system-preset/99-mkosi.preset
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/system/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup@.service.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/systemd/system/systemd-cryptsetup@.service.d/credential.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/udev/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/udev/rules.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/udev/rules.d/10-mkosi-initrd-dm.rules
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-initrd/mkosi.extra/usr/lib/udev/rules.d/10-mkosi-initrd-md.rules
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-arch.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos-fedora/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos-fedora/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos-fedora/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos-fedora/mkosi.conf.d/10-uefi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos/mkosi.conf.d/10-epel.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-centos/mkosi.conf.d/20-epel-packages.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/grub.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/systemd-boot.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/ubuntu-keyring.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-debian-ubuntu/mkosi.conf.d/virtiofsd.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-fedora/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-fedora/mkosi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-fedora/mkosi.conf.d/
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-fedora/mkosi.conf.d/10-uefi.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.conf.d/10-opensuse.conf
usr/lib/python3.13/site-packages/mkosi/resources/mkosi-tools/mkosi.prepare.chroot
usr/lib/python3.13/site-packages/mkosi/resources/mkosi.md
usr/lib/python3.13/site-packages/mkosi/resources/repart/
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/confext.repart.d/
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/confext.repart.d/10-root.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/confext.repart.d/20-root-verity.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/confext.repart.d/30-root-verity-sig.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/portable.repart.d/
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/portable.repart.d/10-root.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/portable.repart.d/20-root-verity.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/portable.repart.d/30-root-verity-sig.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/sysext.repart.d/
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/sysext.repart.d/10-root.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/sysext.repart.d/20-root-verity.conf
usr/lib/python3.13/site-packages/mkosi/resources/repart/definitions/sysext.repart.d/30-root-verity-sig.conf
usr/lib/python3.13/site-packages/mkosi/run.py
usr/lib/python3.13/site-packages/mkosi/sandbox.py
usr/lib/python3.13/site-packages/mkosi/tree.py
usr/lib/python3.13/site-packages/mkosi/types.py
usr/lib/python3.13/site-packages/mkosi/user.py
usr/lib/python3.13/site-packages/mkosi/util.py
usr/lib/python3.13/site-packages/mkosi/versioncomp.py
usr/lib/python3.13/site-packages/mkosi/vmspawn.py
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mkosi.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/mkosi
usr/share/doc/
usr/share/doc/mkosi/
usr/share/doc/mkosi/NEWS.md
usr/share/doc/mkosi/README.md
usr/share/doc/mkosi/bootable.md
usr/share/doc/mkosi/building-rpms-from-source.md
usr/share/doc/mkosi/distribution-policy.md
usr/share/doc/mkosi/index.md
usr/share/doc/mkosi/initrd.md
usr/share/doc/mkosi/sysext.md
usr/share/fish/
usr/share/fish/completions/
usr/share/fish/completions/mkosi.fish
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mkosi.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_mkosi
