# extracted/cni-plugins-1.6.1-1/files
%FILES%
etc/
etc/cni/
etc/cni/net.d/
opt/
opt/cni/
opt/cni/bin/
opt/cni/bin/bandwidth
opt/cni/bin/bridge
opt/cni/bin/dhcp
opt/cni/bin/dummy
opt/cni/bin/firewall
opt/cni/bin/host-device
opt/cni/bin/host-local
opt/cni/bin/ipvlan
opt/cni/bin/loopback
opt/cni/bin/macvlan
opt/cni/bin/portmap
opt/cni/bin/ptp
opt/cni/bin/sbr
opt/cni/bin/static
opt/cni/bin/tap
opt/cni/bin/tuning
opt/cni/bin/vlan
opt/cni/bin/vrf
usr/
usr/lib/
usr/lib/cni/
usr/lib/cni/bandwidth
usr/lib/cni/bridge
usr/lib/cni/dhcp
usr/lib/cni/dummy
usr/lib/cni/firewall
usr/lib/cni/host-device
usr/lib/cni/host-local
usr/lib/cni/ipvlan
usr/lib/cni/loopback
usr/lib/cni/macvlan
usr/lib/cni/portmap
usr/lib/cni/ptp
usr/lib/cni/sbr
usr/lib/cni/static
usr/lib/cni/tap
usr/lib/cni/tuning
usr/lib/cni/vlan
usr/lib/cni/vrf
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cni-dhcp.service
usr/lib/systemd/system/cni-dhcp.socket
usr/share/
usr/share/licenses/
usr/share/licenses/cni-plugins/
usr/share/licenses/cni-plugins/LICENSE
