# extracted/badvpn-1.999.130-7/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/badvpn-ncd
usr/
usr/bin/
usr/bin/badvpn-client
usr/bin/badvpn-flooder
usr/bin/badvpn-ncd
usr/bin/badvpn-ncd-request
usr/bin/badvpn-server
usr/bin/badvpn-tun2socks
usr/bin/badvpn-tunctl
usr/bin/badvpn-udpgw
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/badvpn-ncd.service
usr/lib/systemd/system/badvpn-udpgw.service
usr/share/
usr/share/man/
usr/share/man/man7/
usr/share/man/man7/badvpn.7.gz
usr/share/man/man8/
usr/share/man/man8/badvpn-client.8.gz
usr/share/man/man8/badvpn-server.8.gz
usr/share/man/man8/badvpn-tun2socks.8.gz
