# extracted/plantuml-server-1.2024.6-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/plantuml-server
usr/
usr/bin/
usr/bin/plantuml-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/plantuml-server.service
usr/share/
usr/share/java/
usr/share/java/plantuml-server/
usr/share/java/plantuml-server/dependency/
usr/share/java/plantuml-server/dependency/jetty-runner.jar
usr/share/java/plantuml-server/plantuml.war
usr/share/licenses/
usr/share/licenses/plantuml-server/
usr/share/licenses/plantuml-server/COPYING
