# extracted/cjdns-22-1/files
%FILES%
usr/
usr/bin/
usr/bin/cjdroute
usr/lib/
usr/lib/cjdns/
usr/lib/cjdns/cexec
usr/lib/cjdns/cjdnslog
usr/lib/cjdns/cjdnsshell
usr/lib/cjdns/dumpLinks
usr/lib/cjdns/dumpRumorMill
usr/lib/cjdns/dumptable
usr/lib/cjdns/lib/
usr/lib/cjdns/lib/Semaphore.js
usr/lib/cjdns/lib/cjdnsadmin/
usr/lib/cjdns/lib/cjdnsadmin/bencode.js
usr/lib/cjdns/lib/cjdnsadmin/cjdnsadmin.js
usr/lib/cjdns/lib/makesim.js
usr/lib/cjdns/lib/publicToIp6.js
usr/lib/cjdns/node_modules/
usr/lib/cjdns/node_modules/nthen/
usr/lib/cjdns/node_modules/nthen/.flowconfig
usr/lib/cjdns/node_modules/nthen/.jshintignore
usr/lib/cjdns/node_modules/nthen/.jshintrc
usr/lib/cjdns/node_modules/nthen/.travis.yml
usr/lib/cjdns/node_modules/nthen/index.js
usr/lib/cjdns/node_modules/nthen/package.json
usr/lib/cjdns/node_modules/nthen/readme.md
usr/lib/cjdns/node_modules/nthen/test/
usr/lib/cjdns/node_modules/nthen/test/nthen_test.js
usr/lib/cjdns/node_modules/saferphore/
usr/lib/cjdns/node_modules/saferphore/.flowconfig
usr/lib/cjdns/node_modules/saferphore/.jshintignore
usr/lib/cjdns/node_modules/saferphore/.jshintrc
usr/lib/cjdns/node_modules/saferphore/.travis.yml
usr/lib/cjdns/node_modules/saferphore/index.js
usr/lib/cjdns/node_modules/saferphore/package.json
usr/lib/cjdns/node_modules/saferphore/readme.md
usr/lib/cjdns/node_modules/saferphore/test.js
usr/lib/cjdns/pathfinderTree
usr/lib/cjdns/peerStats
usr/lib/cjdns/ping
usr/lib/cjdns/pingAll
usr/lib/cjdns/search
usr/lib/cjdns/sessionStats
usr/lib/cjdns/splice
usr/lib/cjdns/traceroute
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/cjdns-resume.service
usr/lib/systemd/system/cjdns.service
usr/share/
usr/share/doc/
usr/share/doc/cjdns/
usr/share/doc/cjdns/README.md
usr/share/doc/cjdns/admin-api.md
usr/share/doc/cjdns/configure.md
usr/share/doc/cjdns/djc_layer_model.md
usr/share/doc/cjdns/nat-gateway.md
usr/share/doc/cjdns/network-services.md
usr/share/doc/cjdns/non-root-user.md
usr/share/doc/cjdns/security_specification.md
usr/share/doc/cjdns/shorewall_and_vpn_gateway_howto.md
usr/share/doc/cjdns/tunnel.md
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/cjdroute.conf.5.gz
