# extracted/goaurrpc-1.2.3-2/files
%FILES%
usr/
usr/bin/
usr/bin/goaurrpc
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/goaurrpc.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/goaurrpc.conf
usr/share/
usr/share/doc/
usr/share/doc/goaurrpc/
usr/share/doc/goaurrpc/fromfile.conf.sample
usr/share/doc/goaurrpc/sample.conf.sample
usr/share/licenses/
usr/share/licenses/goaurrpc/
usr/share/licenses/goaurrpc/LICENSE
