# extracted/autorandr-1.15-1/files
%FILES%
etc/
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/autorandr-kde.desktop
etc/xdg/autostart/autorandr.desktop
usr/
usr/bin/
usr/bin/autorandr
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/autorandr-lid-listener.service
usr/lib/systemd/system/autorandr.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/40-monitor-hotplug.rules
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/autorandr
usr/share/doc/
usr/share/doc/autorandr/
usr/share/doc/autorandr/README.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/autorandr.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_autorandr
