# extracted/shadowsocks-rust-1.21.2-1/files
%FILES%
etc/
etc/shadowsocks-rust/
etc/shadowsocks-rust/config_ext_rust.json.example
etc/shadowsocks-rust/config_rust.json.example
usr/
usr/bin/
usr/bin/sslocal
usr/bin/ssmanager
usr/bin/ssserver
usr/bin/ssservice
usr/bin/ssurl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shadowsocks-rust-server@.service
usr/lib/systemd/system/shadowsocks-rust@.service
usr/share/
usr/share/licenses/
usr/share/licenses/shadowsocks-rust/
usr/share/licenses/shadowsocks-rust/LICENSE
