# extracted/memcached-1.6.21-2/files
%FILES%
usr/
usr/bin/
usr/bin/memcached
usr/bin/memcached-tool
usr/include/
usr/include/memcached/
usr/include/memcached/protocol_binary.h
usr/include/memcached/xxhash.h
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/memcached.service
usr/lib/systemd/system/memcached@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/memcached.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/memcached.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/memcached-tool.1.gz
usr/share/man/man1/memcached.1.gz
