# extracted/opensmtpd-7.6.0p1-2/files
%FILES%
etc/
etc/smtpd/
etc/smtpd/aliases
etc/smtpd/smtpd.conf
usr/
usr/bin/
usr/bin/mailq
usr/bin/makemap
usr/bin/newaliases
usr/bin/sendmail
usr/bin/smtp
usr/bin/smtpctl
usr/bin/smtpd
usr/lib/
usr/lib/smtpd/
usr/lib/smtpd/opensmtpd/
usr/lib/smtpd/opensmtpd/encrypt
usr/lib/smtpd/opensmtpd/lockspool
usr/lib/smtpd/opensmtpd/mail.lmtp
usr/lib/smtpd/opensmtpd/mail.local
usr/lib/smtpd/opensmtpd/mail.maildir
usr/lib/smtpd/opensmtpd/mail.mboxfile
usr/lib/smtpd/opensmtpd/mail.mda
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/smtpd.service
usr/lib/systemd/system/smtpd.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/opensmtpd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/opensmtpd/
usr/share/licenses/opensmtpd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/lockspool.1.gz
usr/share/man/man1/smtp.1.gz
usr/share/man/man5/
usr/share/man/man5/aliases.5.gz
usr/share/man/man5/forward.5.gz
usr/share/man/man5/smtpd.conf.5.gz
usr/share/man/man5/table.5.gz
usr/share/man/man7/
usr/share/man/man7/smtpd-filters.7.gz
usr/share/man/man7/smtpd-tables.7.gz
usr/share/man/man8/
usr/share/man/man8/mail.lmtp.8.gz
usr/share/man/man8/mail.local.8.gz
usr/share/man/man8/mail.maildir.8.gz
usr/share/man/man8/mail.mboxfile.8.gz
usr/share/man/man8/mail.mda.8.gz
usr/share/man/man8/makemap.8.gz
usr/share/man/man8/newaliases.8.gz
usr/share/man/man8/sendmail.8.gz
usr/share/man/man8/smtpctl.8.gz
usr/share/man/man8/smtpd.8.gz
