# extracted/unrealircd-6.1.9.1-1/files
%FILES%
etc/
etc/unrealircd/
etc/unrealircd/aliases/
etc/unrealircd/aliases/aliases.conf
etc/unrealircd/aliases/anope.conf
etc/unrealircd/aliases/atheme.conf
etc/unrealircd/aliases/auspice.conf
etc/unrealircd/aliases/cygnus.conf
etc/unrealircd/aliases/epona.conf
etc/unrealircd/aliases/generic.conf
etc/unrealircd/aliases/genericstats.conf
etc/unrealircd/aliases/ircservices.conf
etc/unrealircd/aliases/operstats.conf
etc/unrealircd/badwords.conf
etc/unrealircd/dccallow.conf
etc/unrealircd/examples/
etc/unrealircd/examples/example.conf
etc/unrealircd/examples/example.es.conf
etc/unrealircd/examples/example.fr.conf
etc/unrealircd/examples/example.pt.conf
etc/unrealircd/examples/example.tr.conf
etc/unrealircd/help/
etc/unrealircd/help/help.conf
etc/unrealircd/help/help.de.conf
etc/unrealircd/help/help.es.conf
etc/unrealircd/help/help.fr.conf
etc/unrealircd/help/help.it.conf
etc/unrealircd/help/help.nl.conf
etc/unrealircd/help/help.pl.conf
etc/unrealircd/help/help.ru.conf
etc/unrealircd/help/help.tr.conf
etc/unrealircd/modules.default.conf
etc/unrealircd/modules.optional.conf
etc/unrealircd/modules.sources.list
etc/unrealircd/operclass.default.conf
etc/unrealircd/rpc-class.default.conf
etc/unrealircd/rpc.modules.default.conf
etc/unrealircd/snomasks.default.conf
etc/unrealircd/spamfilter.conf
etc/unrealircd/tls/
etc/unrealircd/tls/curl-ca-bundle.crt
etc/unrealircd/unrealircd
etc/unrealircd/unrealircd.conf
usr/
usr/bin/
usr/bin/unrealircd
usr/bin/unrealircd-upgrade-script
usr/bin/unrealircdctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/unrealircd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/unrealircd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/unrealircd.conf
usr/lib/unrealircd/
usr/lib/unrealircd/account-notify.so
usr/lib/unrealircd/account-tag.so
usr/lib/unrealircd/addmotd.so
usr/lib/unrealircd/addomotd.so
usr/lib/unrealircd/admin.so
usr/lib/unrealircd/antimixedutf8.so
usr/lib/unrealircd/antirandom.so
usr/lib/unrealircd/authprompt.so
usr/lib/unrealircd/away.so
usr/lib/unrealircd/batch.so
usr/lib/unrealircd/blacklist.so
usr/lib/unrealircd/bot-tag.so
usr/lib/unrealircd/botmotd.so
usr/lib/unrealircd/cap.so
usr/lib/unrealircd/central-api.so
usr/lib/unrealircd/central-blocklist.so
usr/lib/unrealircd/certfp.so
usr/lib/unrealircd/chanmodes/
usr/lib/unrealircd/chanmodes/censor.so
usr/lib/unrealircd/chanmodes/chanadmin.so
usr/lib/unrealircd/chanmodes/chanop.so
usr/lib/unrealircd/chanmodes/chanowner.so
usr/lib/unrealircd/chanmodes/delayjoin.so
usr/lib/unrealircd/chanmodes/floodprot.so
usr/lib/unrealircd/chanmodes/halfop.so
usr/lib/unrealircd/chanmodes/history.so
usr/lib/unrealircd/chanmodes/inviteonly.so
usr/lib/unrealircd/chanmodes/isregistered.so
usr/lib/unrealircd/chanmodes/issecure.so
usr/lib/unrealircd/chanmodes/key.so
usr/lib/unrealircd/chanmodes/limit.so
usr/lib/unrealircd/chanmodes/link.so
usr/lib/unrealircd/chanmodes/moderated.so
usr/lib/unrealircd/chanmodes/nocolor.so
usr/lib/unrealircd/chanmodes/noctcp.so
usr/lib/unrealircd/chanmodes/noexternalmsgs.so
usr/lib/unrealircd/chanmodes/noinvite.so
usr/lib/unrealircd/chanmodes/nokick.so
usr/lib/unrealircd/chanmodes/noknock.so
usr/lib/unrealircd/chanmodes/nonickchange.so
usr/lib/unrealircd/chanmodes/nonotice.so
usr/lib/unrealircd/chanmodes/operonly.so
usr/lib/unrealircd/chanmodes/permanent.so
usr/lib/unrealircd/chanmodes/private.so
usr/lib/unrealircd/chanmodes/regonly.so
usr/lib/unrealircd/chanmodes/regonlyspeak.so
usr/lib/unrealircd/chanmodes/secret.so
usr/lib/unrealircd/chanmodes/secureonly.so
usr/lib/unrealircd/chanmodes/stripcolor.so
usr/lib/unrealircd/chanmodes/topiclimit.so
usr/lib/unrealircd/chanmodes/voice.so
usr/lib/unrealircd/channel-context.so
usr/lib/unrealircd/channeldb.so
usr/lib/unrealircd/charsys.so
usr/lib/unrealircd/chathistory.so
usr/lib/unrealircd/chghost.so
usr/lib/unrealircd/chgident.so
usr/lib/unrealircd/chgname.so
usr/lib/unrealircd/clienttagdeny.so
usr/lib/unrealircd/cloak_md5.so
usr/lib/unrealircd/cloak_none.so
usr/lib/unrealircd/cloak_sha256.so
usr/lib/unrealircd/close.so
usr/lib/unrealircd/connect-flood.so
usr/lib/unrealircd/connect.so
usr/lib/unrealircd/connthrottle.so
usr/lib/unrealircd/creationtime.so
usr/lib/unrealircd/crule.so
usr/lib/unrealircd/cycle.so
usr/lib/unrealircd/dccallow.so
usr/lib/unrealircd/dccdeny.so
usr/lib/unrealircd/echo-message.so
usr/lib/unrealircd/eos.so
usr/lib/unrealircd/extbans/
usr/lib/unrealircd/extbans/account.so
usr/lib/unrealircd/extbans/asn.so
usr/lib/unrealircd/extbans/certfp.so
usr/lib/unrealircd/extbans/country.so
usr/lib/unrealircd/extbans/flood.so
usr/lib/unrealircd/extbans/inchannel.so
usr/lib/unrealircd/extbans/inherit.so
usr/lib/unrealircd/extbans/join.so
usr/lib/unrealircd/extbans/msgbypass.so
usr/lib/unrealircd/extbans/nickchange.so
usr/lib/unrealircd/extbans/operclass.so
usr/lib/unrealircd/extbans/partmsg.so
usr/lib/unrealircd/extbans/quiet.so
usr/lib/unrealircd/extbans/realname.so
usr/lib/unrealircd/extbans/securitygroup.so
usr/lib/unrealircd/extbans/textban.so
usr/lib/unrealircd/extbans/timedban.so
usr/lib/unrealircd/extended-monitor.so
usr/lib/unrealircd/extjwt.so
usr/lib/unrealircd/geoip-tag.so
usr/lib/unrealircd/geoip_base.so
usr/lib/unrealircd/geoip_csv.so
usr/lib/unrealircd/globops.so
usr/lib/unrealircd/help.so
usr/lib/unrealircd/hideserver.so
usr/lib/unrealircd/history.so
usr/lib/unrealircd/history_backend_mem.so
usr/lib/unrealircd/history_backend_null.so
usr/lib/unrealircd/ident_lookup.so
usr/lib/unrealircd/invite.so
usr/lib/unrealircd/ircops.so
usr/lib/unrealircd/ison.so
usr/lib/unrealircd/issued-by-tag.so
usr/lib/unrealircd/join.so
usr/lib/unrealircd/jointhrottle.so
usr/lib/unrealircd/json-log-tag.so
usr/lib/unrealircd/jumpserver.so
usr/lib/unrealircd/kick.so
usr/lib/unrealircd/kill.so
usr/lib/unrealircd/knock.so
usr/lib/unrealircd/labeled-response.so
usr/lib/unrealircd/lag.so
usr/lib/unrealircd/link-security.so
usr/lib/unrealircd/links.so
usr/lib/unrealircd/list.so
usr/lib/unrealircd/locops.so
usr/lib/unrealircd/lusers.so
usr/lib/unrealircd/map.so
usr/lib/unrealircd/max-unknown-connections-per-ip.so
usr/lib/unrealircd/md.so
usr/lib/unrealircd/message-ids.so
usr/lib/unrealircd/message-tags.so
usr/lib/unrealircd/message.so
usr/lib/unrealircd/mkpasswd.so
usr/lib/unrealircd/mode.so
usr/lib/unrealircd/monitor.so
usr/lib/unrealircd/motd.so
usr/lib/unrealircd/names.so
usr/lib/unrealircd/netinfo.so
usr/lib/unrealircd/nick.so
usr/lib/unrealircd/no-implicit-names.so
usr/lib/unrealircd/nocodes.so
usr/lib/unrealircd/oper.so
usr/lib/unrealircd/operinfo.so
usr/lib/unrealircd/opermotd.so
usr/lib/unrealircd/part.so
usr/lib/unrealircd/pass.so
usr/lib/unrealircd/pingpong.so
usr/lib/unrealircd/plaintext-policy.so
usr/lib/unrealircd/protoctl.so
usr/lib/unrealircd/quit.so
usr/lib/unrealircd/real-quit-reason.so
usr/lib/unrealircd/reply-tag.so
usr/lib/unrealircd/reputation.so
usr/lib/unrealircd/require-module.so
usr/lib/unrealircd/restrict-commands.so
usr/lib/unrealircd/rmtkl.so
usr/lib/unrealircd/rpc/
usr/lib/unrealircd/rpc/channel.so
usr/lib/unrealircd/rpc/log.so
usr/lib/unrealircd/rpc/name_ban.so
usr/lib/unrealircd/rpc/rpc.so
usr/lib/unrealircd/rpc/server.so
usr/lib/unrealircd/rpc/server_ban.so
usr/lib/unrealircd/rpc/server_ban_exception.so
usr/lib/unrealircd/rpc/spamfilter.so
usr/lib/unrealircd/rpc/stats.so
usr/lib/unrealircd/rpc/user.so
usr/lib/unrealircd/rpc/whowas.so
usr/lib/unrealircd/rules.so
usr/lib/unrealircd/sajoin.so
usr/lib/unrealircd/samode.so
usr/lib/unrealircd/sapart.so
usr/lib/unrealircd/sasl.so
usr/lib/unrealircd/sdesc.so
usr/lib/unrealircd/sendsno.so
usr/lib/unrealircd/sendumode.so
usr/lib/unrealircd/server-time.so
usr/lib/unrealircd/server.so
usr/lib/unrealircd/sethost.so
usr/lib/unrealircd/setident.so
usr/lib/unrealircd/setname.so
usr/lib/unrealircd/silence.so
usr/lib/unrealircd/sinfo.so
usr/lib/unrealircd/sjoin.so
usr/lib/unrealircd/slog.so
usr/lib/unrealircd/spamreport.so
usr/lib/unrealircd/sqline.so
usr/lib/unrealircd/squit.so
usr/lib/unrealircd/sreply.so
usr/lib/unrealircd/staff.so
usr/lib/unrealircd/standard-replies.so
usr/lib/unrealircd/starttls.so
usr/lib/unrealircd/stats.so
usr/lib/unrealircd/sts.so
usr/lib/unrealircd/svsjoin.so
usr/lib/unrealircd/svskill.so
usr/lib/unrealircd/svslogin.so
usr/lib/unrealircd/svslusers.so
usr/lib/unrealircd/svsmode.so
usr/lib/unrealircd/svsmotd.so
usr/lib/unrealircd/svsnick.so
usr/lib/unrealircd/svsnline.so
usr/lib/unrealircd/svsnolag.so
usr/lib/unrealircd/svsnoop.so
usr/lib/unrealircd/svso.so
usr/lib/unrealircd/svspart.so
usr/lib/unrealircd/svssilence.so
usr/lib/unrealircd/svssno.so
usr/lib/unrealircd/svswatch.so
usr/lib/unrealircd/swhois.so
usr/lib/unrealircd/targetfloodprot.so
usr/lib/unrealircd/third/
usr/lib/unrealircd/time.so
usr/lib/unrealircd/tkl.so
usr/lib/unrealircd/tkldb.so
usr/lib/unrealircd/tline.so
usr/lib/unrealircd/tls_antidos.so
usr/lib/unrealircd/tls_cipher.so
usr/lib/unrealircd/topic.so
usr/lib/unrealircd/trace.so
usr/lib/unrealircd/tsctl.so
usr/lib/unrealircd/typing-indicator.so
usr/lib/unrealircd/umode2.so
usr/lib/unrealircd/unreal_server_compat.so
usr/lib/unrealircd/unsqline.so
usr/lib/unrealircd/user.so
usr/lib/unrealircd/userhost-tag.so
usr/lib/unrealircd/userhost.so
usr/lib/unrealircd/userip-tag.so
usr/lib/unrealircd/userip.so
usr/lib/unrealircd/usermodes/
usr/lib/unrealircd/usermodes/bot.so
usr/lib/unrealircd/usermodes/censor.so
usr/lib/unrealircd/usermodes/noctcp.so
usr/lib/unrealircd/usermodes/nokick.so
usr/lib/unrealircd/usermodes/privacy.so
usr/lib/unrealircd/usermodes/privdeaf.so
usr/lib/unrealircd/usermodes/regonlymsg.so
usr/lib/unrealircd/usermodes/secureonlymsg.so
usr/lib/unrealircd/usermodes/servicebot.so
usr/lib/unrealircd/usermodes/showwhois.so
usr/lib/unrealircd/usermodes/wallops.so
usr/lib/unrealircd/vhost.so
usr/lib/unrealircd/watch-backend.so
usr/lib/unrealircd/watch.so
usr/lib/unrealircd/webirc.so
usr/lib/unrealircd/webredir.so
usr/lib/unrealircd/webserver.so
usr/lib/unrealircd/websocket.so
usr/lib/unrealircd/websocket_common.so
usr/lib/unrealircd/who_old.so
usr/lib/unrealircd/whois.so
usr/lib/unrealircd/whowas.so
usr/lib/unrealircd/whowasdb.so
usr/lib/unrealircd/whox.so
usr/share/
usr/share/doc/
usr/share/doc/unrealircd/
usr/share/doc/unrealircd/Authors
usr/share/doc/unrealircd/KEYS
usr/share/doc/unrealircd/RELEASE-NOTES.md
usr/share/doc/unrealircd/coding-guidelines
usr/share/doc/unrealircd/tao.of.irc
var/
var/cache/
var/cache/unrealircd/
var/lib/
var/lib/unrealircd/
var/log/
var/log/unrealircd/
