# extracted/ipvsadm-1.31-3/files
%FILES%
etc/
etc/ipvsadm.conf
usr/
usr/bin/
usr/bin/ipvsadm
usr/bin/ipvsadm-restore
usr/bin/ipvsadm-save
usr/lib/
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/ipvsadm
usr/lib/systemd/system/
usr/lib/systemd/system/ipvsadm.service
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/ipvsadm-restore.8.gz
usr/share/man/man8/ipvsadm-save.8.gz
usr/share/man/man8/ipvsadm.8.gz
