# extracted/dns-over-https-2.3.7-1/files
%FILES%
etc/
etc/NetworkManager/
etc/NetworkManager/dispatcher.d/
etc/NetworkManager/dispatcher.d/doh-client
etc/NetworkManager/dispatcher.d/doh-server
etc/dns-over-https/
etc/dns-over-https/doh-client.conf
etc/dns-over-https/doh-server.conf
usr/
usr/bin/
usr/bin/doh-client
usr/bin/doh-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/doh-client.service
usr/lib/systemd/system/doh-server.service
usr/share/
usr/share/licenses/
usr/share/licenses/dns-over-https/
usr/share/licenses/dns-over-https/LICENSE
