# extracted/sniproxy-0.6.1-2/files
%FILES%
etc/
etc/sniproxy.conf
usr/
usr/bin/
usr/bin/sniproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sniproxy.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sniproxy.conf
usr/share/
usr/share/doc/
usr/share/doc/sniproxy/
usr/share/doc/sniproxy/sniproxy.conf
usr/share/licenses/
usr/share/licenses/sniproxy/
usr/share/licenses/sniproxy/COPYING
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/sniproxy.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/sniproxy.8.gz
