# extracted/pkgstats-3.3.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/pkgstats
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pkgstats.service
usr/lib/systemd/system/pkgstats.timer
usr/lib/systemd/system/timers.target.wants/
usr/lib/systemd/system/timers.target.wants/pkgstats.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/pkgstats
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/pkgstats.fish
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_pkgstats
