# extracted/udpxy-1.0.25.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/udpxrec
usr/bin/udpxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/udpxy.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/udpxrec.1.gz
usr/share/man/man1/udpxy.1.gz
