# extracted/shadowsocks-3.0.0a.20180219-7/files
%FILES%
etc/
etc/shadowsocks/
etc/shadowsocks/example.json
usr/
usr/bin/
usr/bin/sslocal
usr/bin/ssserver
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/PKG-INFO
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/SOURCES.txt
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/dependency_links.txt
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/entry_points.txt
usr/lib/python3.12/site-packages/shadowsocks-3.0.0-py3.12.egg-info/top_level.txt
usr/lib/python3.12/site-packages/shadowsocks/
usr/lib/python3.12/site-packages/shadowsocks/__init__.py
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/asyncdns.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/asyncdns.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/common.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/common.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/cryptor.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/cryptor.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/daemon.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/daemon.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/eventloop.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/eventloop.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/local.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/local.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/lru_cache.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/lru_cache.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/manager.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/manager.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/server.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/server.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/shell.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/shell.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/tcprelay.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/tcprelay.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/tunnel.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/tunnel.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/udprelay.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/__pycache__/udprelay.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/asyncdns.py
usr/lib/python3.12/site-packages/shadowsocks/common.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/
usr/lib/python3.12/site-packages/shadowsocks/crypto/__init__.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/aead.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/aead.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/hkdf.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/hkdf.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/mbedtls.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/mbedtls.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/openssl.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/openssl.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/rc4_md5.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/rc4_md5.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/sodium.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/sodium.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/table.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/table.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/util.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/__pycache__/util.cpython-312.pyc
usr/lib/python3.12/site-packages/shadowsocks/crypto/aead.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/hkdf.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/mbedtls.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/openssl.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/rc4_md5.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/sodium.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/table.py
usr/lib/python3.12/site-packages/shadowsocks/crypto/util.py
usr/lib/python3.12/site-packages/shadowsocks/cryptor.py
usr/lib/python3.12/site-packages/shadowsocks/daemon.py
usr/lib/python3.12/site-packages/shadowsocks/eventloop.py
usr/lib/python3.12/site-packages/shadowsocks/local.py
usr/lib/python3.12/site-packages/shadowsocks/lru_cache.py
usr/lib/python3.12/site-packages/shadowsocks/manager.py
usr/lib/python3.12/site-packages/shadowsocks/server.py
usr/lib/python3.12/site-packages/shadowsocks/shell.py
usr/lib/python3.12/site-packages/shadowsocks/tcprelay.py
usr/lib/python3.12/site-packages/shadowsocks/tunnel.py
usr/lib/python3.12/site-packages/shadowsocks/udprelay.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shadowsocks-server@.service
usr/lib/systemd/system/shadowsocks@.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/sslocal.1.gz
usr/share/man/man1/ssserver.1.gz
