# extracted/pesign-116-2/files
%FILES%
etc/
etc/pesign/
etc/pesign/groups
etc/pesign/users
etc/pki/
etc/pki/pesign/
etc/popt.d/
etc/popt.d/pesign.popt
usr/
usr/bin/
usr/bin/authvar
usr/bin/efikeygen
usr/bin/pesigcheck
usr/bin/pesign
usr/bin/pesign-client
usr/bin/pesum
usr/lib/
usr/lib/pesign/
usr/lib/pesign/pesign-authorize
usr/lib/pesign/pesign-rpmbuild-helper
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pesign-create-db.service
usr/lib/systemd/system/pesign.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pesign.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pesign.conf
usr/share/
usr/share/doc/
usr/share/doc/pesign-116/
usr/share/doc/pesign-116/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/authvar.1.gz
usr/share/man/man1/efikeygen.1.gz
usr/share/man/man1/pesigcheck.1.gz
usr/share/man/man1/pesign-client.1.gz
usr/share/man/man1/pesign.1.gz
