# extracted/nss-pam-ldapd-0.9.12-4/files
%FILES%
etc/
etc/nslcd.conf
usr/
usr/bin/
usr/bin/chsh.ldap
usr/bin/getent.ldap
usr/bin/nslcd
usr/lib/
usr/lib/libnss_ldap.so.2
usr/lib/security/
usr/lib/security/pam_ldap.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nslcd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/nslcd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/nslcd.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/chsh.ldap.1.gz
usr/share/man/man1/getent.ldap.1.gz
usr/share/man/man5/
usr/share/man/man5/nslcd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/nslcd.8.gz
usr/share/man/man8/pam_ldap.8.gz
usr/share/nslcd-utils/
usr/share/nslcd-utils/__pycache__/
usr/share/nslcd-utils/__pycache__/chsh.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/chsh.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/cmdline.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/cmdline.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/constants.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/constants.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/getent.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/getent.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/nslcd.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/nslcd.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/shells.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/shells.cpython-312.pyc
usr/share/nslcd-utils/__pycache__/users.cpython-312.opt-1.pyc
usr/share/nslcd-utils/__pycache__/users.cpython-312.pyc
usr/share/nslcd-utils/chsh.py
usr/share/nslcd-utils/cmdline.py
usr/share/nslcd-utils/constants.py
usr/share/nslcd-utils/getent.py
usr/share/nslcd-utils/nslcd.py
usr/share/nslcd-utils/shells.py
usr/share/nslcd-utils/users.py
