# extracted/unbound-1.22.0-2/files
%FILES%
etc/
etc/unbound/
etc/unbound/unbound.conf
usr/
usr/bin/
usr/bin/unbound
usr/bin/unbound-anchor
usr/bin/unbound-checkconf
usr/bin/unbound-control
usr/bin/unbound-control-setup
usr/bin/unbound-host
usr/include/
usr/include/unbound.h
usr/lib/
usr/lib/libunbound.so
usr/lib/libunbound.so.8
usr/lib/libunbound.so.8.1.30
usr/lib/pkgconfig/
usr/lib/pkgconfig/libunbound.pc
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/_unbound.so
usr/lib/python3.13/site-packages/unbound.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/unbound.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/unbound.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/unbound.conf
usr/share/
usr/share/libalpm/
usr/share/libalpm/hooks/
usr/share/libalpm/hooks/unbound-trusted-key.hook
usr/share/licenses/
usr/share/licenses/unbound/
usr/share/licenses/unbound/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/unbound-host.1.gz
usr/share/man/man3/
usr/share/man/man3/libunbound.3.gz
usr/share/man/man3/ub_cancel.3.gz
usr/share/man/man3/ub_ctx.3.gz
usr/share/man/man3/ub_ctx_add_ta.3.gz
usr/share/man/man3/ub_ctx_add_ta_file.3.gz
usr/share/man/man3/ub_ctx_async.3.gz
usr/share/man/man3/ub_ctx_config.3.gz
usr/share/man/man3/ub_ctx_create.3.gz
usr/share/man/man3/ub_ctx_data_add.3.gz
usr/share/man/man3/ub_ctx_data_remove.3.gz
usr/share/man/man3/ub_ctx_debuglevel.3.gz
usr/share/man/man3/ub_ctx_debugout.3.gz
usr/share/man/man3/ub_ctx_delete.3.gz
usr/share/man/man3/ub_ctx_get_option.3.gz
usr/share/man/man3/ub_ctx_hosts.3.gz
usr/share/man/man3/ub_ctx_print_local_zones.3.gz
usr/share/man/man3/ub_ctx_resolvconf.3.gz
usr/share/man/man3/ub_ctx_set_fwd.3.gz
usr/share/man/man3/ub_ctx_set_option.3.gz
usr/share/man/man3/ub_ctx_trustedkeys.3.gz
usr/share/man/man3/ub_ctx_zone_add.3.gz
usr/share/man/man3/ub_ctx_zone_remove.3.gz
usr/share/man/man3/ub_fd.3.gz
usr/share/man/man3/ub_poll.3.gz
usr/share/man/man3/ub_process.3.gz
usr/share/man/man3/ub_resolve.3.gz
usr/share/man/man3/ub_resolve_async.3.gz
usr/share/man/man3/ub_resolve_free.3.gz
usr/share/man/man3/ub_result.3.gz
usr/share/man/man3/ub_strerror.3.gz
usr/share/man/man3/ub_wait.3.gz
usr/share/man/man5/
usr/share/man/man5/unbound.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/unbound-anchor.8.gz
usr/share/man/man8/unbound-checkconf.8.gz
usr/share/man/man8/unbound-control-setup.8.gz
usr/share/man/man8/unbound-control.8.gz
usr/share/man/man8/unbound.8.gz
