# extracted/kubeadm-1.31.3-1/files
%FILES%
usr/
usr/bin/
usr/bin/kubeadm
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/kubelet.service.d/
usr/lib/systemd/system/kubelet.service.d/10-kubeadm.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/kubeadm
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/kubeadm-alpha.1.gz
usr/share/man/man1/kubeadm-certs-certificate-key.1.gz
usr/share/man/man1/kubeadm-certs-check-expiration.1.gz
usr/share/man/man1/kubeadm-certs-generate-csr.1.gz
usr/share/man/man1/kubeadm-certs-renew-admin.conf.1.gz
usr/share/man/man1/kubeadm-certs-renew-all.1.gz
usr/share/man/man1/kubeadm-certs-renew-apiserver-etcd-client.1.gz
usr/share/man/man1/kubeadm-certs-renew-apiserver-kubelet-client.1.gz
usr/share/man/man1/kubeadm-certs-renew-apiserver.1.gz
usr/share/man/man1/kubeadm-certs-renew-controller-manager.conf.1.gz
usr/share/man/man1/kubeadm-certs-renew-etcd-healthcheck-client.1.gz
usr/share/man/man1/kubeadm-certs-renew-etcd-peer.1.gz
usr/share/man/man1/kubeadm-certs-renew-etcd-server.1.gz
usr/share/man/man1/kubeadm-certs-renew-front-proxy-client.1.gz
usr/share/man/man1/kubeadm-certs-renew-scheduler.conf.1.gz
usr/share/man/man1/kubeadm-certs-renew-super-admin.conf.1.gz
usr/share/man/man1/kubeadm-certs-renew.1.gz
usr/share/man/man1/kubeadm-certs.1.gz
usr/share/man/man1/kubeadm-completion.1.gz
usr/share/man/man1/kubeadm-config-images-list.1.gz
usr/share/man/man1/kubeadm-config-images-pull.1.gz
usr/share/man/man1/kubeadm-config-images.1.gz
usr/share/man/man1/kubeadm-config-migrate.1.gz
usr/share/man/man1/kubeadm-config-print-init-defaults.1.gz
usr/share/man/man1/kubeadm-config-print-join-defaults.1.gz
usr/share/man/man1/kubeadm-config-print-reset-defaults.1.gz
usr/share/man/man1/kubeadm-config-print-upgrade-defaults.1.gz
usr/share/man/man1/kubeadm-config-print.1.gz
usr/share/man/man1/kubeadm-config-validate.1.gz
usr/share/man/man1/kubeadm-config.1.gz
usr/share/man/man1/kubeadm-init-phase-addon-all.1.gz
usr/share/man/man1/kubeadm-init-phase-addon-coredns.1.gz
usr/share/man/man1/kubeadm-init-phase-addon-kube-proxy.1.gz
usr/share/man/man1/kubeadm-init-phase-addon.1.gz
usr/share/man/man1/kubeadm-init-phase-bootstrap-token.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-all.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-apiserver-etcd-client.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-apiserver-kubelet-client.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-apiserver.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-ca.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-etcd-ca.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-etcd-healthcheck-client.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-etcd-peer.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-etcd-server.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-front-proxy-ca.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-front-proxy-client.1.gz
usr/share/man/man1/kubeadm-init-phase-certs-sa.1.gz
usr/share/man/man1/kubeadm-init-phase-certs.1.gz
usr/share/man/man1/kubeadm-init-phase-control-plane-all.1.gz
usr/share/man/man1/kubeadm-init-phase-control-plane-apiserver.1.gz
usr/share/man/man1/kubeadm-init-phase-control-plane-controller-manager.1.gz
usr/share/man/man1/kubeadm-init-phase-control-plane-scheduler.1.gz
usr/share/man/man1/kubeadm-init-phase-control-plane.1.gz
usr/share/man/man1/kubeadm-init-phase-etcd-local.1.gz
usr/share/man/man1/kubeadm-init-phase-etcd.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-admin.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-all.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-controller-manager.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-kubelet.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-scheduler.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig-super-admin.1.gz
usr/share/man/man1/kubeadm-init-phase-kubeconfig.1.gz
usr/share/man/man1/kubeadm-init-phase-kubelet-finalize-all.1.gz
usr/share/man/man1/kubeadm-init-phase-kubelet-finalize-enable-client-cert-rotation.1.gz
usr/share/man/man1/kubeadm-init-phase-kubelet-finalize-experimental-cert-rotation.1.gz
usr/share/man/man1/kubeadm-init-phase-kubelet-finalize.1.gz
usr/share/man/man1/kubeadm-init-phase-kubelet-start.1.gz
usr/share/man/man1/kubeadm-init-phase-mark-control-plane.1.gz
usr/share/man/man1/kubeadm-init-phase-preflight.1.gz
usr/share/man/man1/kubeadm-init-phase-show-join-command.1.gz
usr/share/man/man1/kubeadm-init-phase-upload-certs.1.gz
usr/share/man/man1/kubeadm-init-phase-upload-config-all.1.gz
usr/share/man/man1/kubeadm-init-phase-upload-config-kubeadm.1.gz
usr/share/man/man1/kubeadm-init-phase-upload-config-kubelet.1.gz
usr/share/man/man1/kubeadm-init-phase-upload-config.1.gz
usr/share/man/man1/kubeadm-init-phase.1.gz
usr/share/man/man1/kubeadm-init.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-join-all.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-join-etcd.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-join-mark-control-plane.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-join.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare-all.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare-certs.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare-control-plane.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare-download-certs.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare-kubeconfig.1.gz
usr/share/man/man1/kubeadm-join-phase-control-plane-prepare.1.gz
usr/share/man/man1/kubeadm-join-phase-kubelet-start.1.gz
usr/share/man/man1/kubeadm-join-phase-preflight.1.gz
usr/share/man/man1/kubeadm-join-phase-wait-control-plane.1.gz
usr/share/man/man1/kubeadm-join-phase.1.gz
usr/share/man/man1/kubeadm-join.1.gz
usr/share/man/man1/kubeadm-kubeconfig-user.1.gz
usr/share/man/man1/kubeadm-kubeconfig.1.gz
usr/share/man/man1/kubeadm-reset-phase-cleanup-node.1.gz
usr/share/man/man1/kubeadm-reset-phase-preflight.1.gz
usr/share/man/man1/kubeadm-reset-phase-remove-etcd-member.1.gz
usr/share/man/man1/kubeadm-reset-phase.1.gz
usr/share/man/man1/kubeadm-reset.1.gz
usr/share/man/man1/kubeadm-token-create.1.gz
usr/share/man/man1/kubeadm-token-delete.1.gz
usr/share/man/man1/kubeadm-token-generate.1.gz
usr/share/man/man1/kubeadm-token-list.1.gz
usr/share/man/man1/kubeadm-token.1.gz
usr/share/man/man1/kubeadm-upgrade-apply.1.gz
usr/share/man/man1/kubeadm-upgrade-diff.1.gz
usr/share/man/man1/kubeadm-upgrade-node-phase-control-plane.1.gz
usr/share/man/man1/kubeadm-upgrade-node-phase-kubelet-config.1.gz
usr/share/man/man1/kubeadm-upgrade-node-phase-preflight.1.gz
usr/share/man/man1/kubeadm-upgrade-node-phase.1.gz
usr/share/man/man1/kubeadm-upgrade-node.1.gz
usr/share/man/man1/kubeadm-upgrade-plan.1.gz
usr/share/man/man1/kubeadm-upgrade.1.gz
usr/share/man/man1/kubeadm-version.1.gz
usr/share/man/man1/kubeadm.1.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_kubeadm
