# extracted/bcachefs-tools-3:1.13.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/bcachefs
usr/bin/fsck.bcachefs
usr/bin/fsck.fuse.bcachefs
usr/bin/mkfs.bcachefs
usr/bin/mkfs.fuse.bcachefs
usr/bin/mount.bcachefs
usr/bin/mount.fuse.bcachefs
usr/lib/
usr/lib/bcachefsck_all
usr/lib/bcachefsck_fail
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/bcachefsck@.service
usr/lib/systemd/system/bcachefsck_all.service
usr/lib/systemd/system/bcachefsck_all.timer
usr/lib/systemd/system/bcachefsck_all_fail.service
usr/lib/systemd/system/bcachefsck_fail@.service
usr/lib/systemd/system/system-bcachefsck.slice
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/64-bcachefs.rules
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/bcachefs
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/bcachefs.fish
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/bcachefs.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_bcachefs
