# extracted/virtualbox-guest-utils-nox-7.1.4-2/files
%FILES%
usr/
usr/bin/
usr/bin/VBoxControl
usr/bin/VBoxService
usr/lib/
usr/lib/security/
usr/lib/security/pam_vbox.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/vboxservice.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/virtualbox-guest-utils.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-vboxguest.rules
usr/share/
usr/share/licenses/
usr/share/licenses/virtualbox-guest-utils-nox/
usr/share/licenses/virtualbox-guest-utils-nox/LICENSE
usr/share/licenses/virtualbox-guest-utils-nox/LICENSE.CDDL
