# extracted/spampd-2.61-2/files
%FILES%
usr/
usr/bin/
usr/bin/spampd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/spampd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/spampd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/spampd.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/spampd.8.gz
