# extracted/fastd-22-4/files
%FILES%
etc/
etc/fastd/
usr/
usr/bin/
usr/bin/fastd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/fastd@.service
usr/share/
usr/share/licenses/
usr/share/licenses/fastd/
usr/share/licenses/fastd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/fastd.1.gz
