# extracted/munge-0.5.16-1/files
%FILES%
etc/
etc/default/
etc/default/munge
etc/munge/
usr/
usr/bin/
usr/bin/munge
usr/bin/munged
usr/bin/mungekey
usr/bin/remunge
usr/bin/unmunge
usr/include/
usr/include/munge.h
usr/lib/
usr/lib/libmunge.so
usr/lib/libmunge.so.2
usr/lib/libmunge.so.2.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/munge.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/munge.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/munge.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/munge.conf
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/munge.1.gz
usr/share/man/man1/remunge.1.gz
usr/share/man/man1/unmunge.1.gz
usr/share/man/man3/
usr/share/man/man3/munge.3.gz
usr/share/man/man3/munge_ctx.3.gz
usr/share/man/man3/munge_ctx_copy.3.gz
usr/share/man/man3/munge_ctx_create.3.gz
usr/share/man/man3/munge_ctx_destroy.3.gz
usr/share/man/man3/munge_ctx_get.3.gz
usr/share/man/man3/munge_ctx_set.3.gz
usr/share/man/man3/munge_ctx_strerror.3.gz
usr/share/man/man3/munge_decode.3.gz
usr/share/man/man3/munge_encode.3.gz
usr/share/man/man3/munge_enum.3.gz
usr/share/man/man3/munge_enum_int_to_str.3.gz
usr/share/man/man3/munge_enum_is_valid.3.gz
usr/share/man/man3/munge_enum_str_to_int.3.gz
usr/share/man/man3/munge_strerror.3.gz
usr/share/man/man7/
usr/share/man/man7/munge.7.gz
usr/share/man/man8/
usr/share/man/man8/munged.8.gz
usr/share/man/man8/mungekey.8.gz
var/
var/lib/
var/lib/munge/
var/log/
var/log/munge/
