# extracted/php-fpm-8.3.14-1/files
%FILES%
etc/
etc/php/
etc/php/php-fpm.conf
etc/php/php-fpm.d/
etc/php/php-fpm.d/www.conf
usr/
usr/bin/
usr/bin/php-fpm
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/php-fpm.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/php-fpm.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/php-fpm.8.gz
usr/share/php/
usr/share/php/fpm/
usr/share/php/fpm/status.html
