# extracted/clash-1.18.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/clash
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/clash@.service
usr/lib/systemd/user/
usr/lib/systemd/user/clash.service
usr/share/
usr/share/licenses/
usr/share/licenses/clash/
usr/share/licenses/clash/LICENSE
