# extracted/arpwatch-3.6-1/files
%FILES%
usr/
usr/bin/
usr/bin/arpsnmp
usr/bin/arpwatch
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/arpwatch@.service
usr/share/
usr/share/licenses/
usr/share/licenses/arpwatch/
usr/share/licenses/arpwatch/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/arpsnmp.8.gz
usr/share/man/man8/arpwatch.8.gz
var/
var/lib/
var/lib/arpwatch/
var/lib/arpwatch/ethercodes.dat
