# extracted/sssd-2.10.1-4/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/sssd
etc/pam.d/
etc/pam.d/sssd-shadowutils
etc/sssd/
etc/sssd/conf.d/
etc/sssd/pki/
usr/
usr/bin/
usr/bin/sss_cache
usr/bin/sss_debuglevel
usr/bin/sss_obfuscate
usr/bin/sss_override
usr/bin/sss_seed
usr/bin/sss_ssh_authorizedkeys
usr/bin/sss_ssh_knownhosts
usr/bin/sss_ssh_knownhostsproxy
usr/bin/sssctl
usr/bin/sssd
usr/include/
usr/include/ipa_hbac.h
usr/include/sss_certmap.h
usr/include/sss_idmap.h
usr/include/sss_nss_idmap.h
usr/lib/
usr/lib/cifs-utils/
usr/lib/cifs-utils/cifs_idmap_sss.so
usr/lib/krb5/
usr/lib/krb5/plugins/
usr/lib/krb5/plugins/authdata/
usr/lib/krb5/plugins/authdata/sssd_pac_plugin.so
usr/lib/krb5/plugins/libkrb5/
usr/lib/krb5/plugins/libkrb5/sssd_krb5_locator_plugin.so
usr/lib/libipa_hbac.so
usr/lib/libipa_hbac.so.0
usr/lib/libipa_hbac.so.0.1.0
usr/lib/libnfsidmap/
usr/lib/libnfsidmap/sss.so
usr/lib/libnss_sss.so.2
usr/lib/libsss_certmap.so
usr/lib/libsss_certmap.so.0
usr/lib/libsss_certmap.so.0.2.0
usr/lib/libsss_idmap.so
usr/lib/libsss_idmap.so.0
usr/lib/libsss_idmap.so.0.5.1
usr/lib/libsss_nss_idmap.so
usr/lib/libsss_nss_idmap.so.0
usr/lib/libsss_nss_idmap.so.0.6.0
usr/lib/libsss_sudo.so
usr/lib/libsubid_sss.so
usr/lib/pkgconfig/
usr/lib/pkgconfig/ipa_hbac.pc
usr/lib/pkgconfig/sss_certmap.pc
usr/lib/pkgconfig/sss_idmap.pc
usr/lib/pkgconfig/sss_nss_idmap.pc
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/SSSDConfig-2.10.1-py3.13.egg-info/
usr/lib/python3.13/site-packages/SSSDConfig-2.10.1-py3.13.egg-info/PKG-INFO
usr/lib/python3.13/site-packages/SSSDConfig-2.10.1-py3.13.egg-info/SOURCES.txt
usr/lib/python3.13/site-packages/SSSDConfig-2.10.1-py3.13.egg-info/dependency_links.txt
usr/lib/python3.13/site-packages/SSSDConfig-2.10.1-py3.13.egg-info/top_level.txt
usr/lib/python3.13/site-packages/SSSDConfig/
usr/lib/python3.13/site-packages/SSSDConfig/__init__.py
usr/lib/python3.13/site-packages/SSSDConfig/__pycache__/
usr/lib/python3.13/site-packages/SSSDConfig/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/SSSDConfig/__pycache__/ipachangeconf.cpython-313.pyc
usr/lib/python3.13/site-packages/SSSDConfig/__pycache__/sssdoptions.cpython-313.pyc
usr/lib/python3.13/site-packages/SSSDConfig/ipachangeconf.py
usr/lib/python3.13/site-packages/SSSDConfig/sssdoptions.py
usr/lib/python3.13/site-packages/pyhbac.so
usr/lib/python3.13/site-packages/pysss.so
usr/lib/python3.13/site-packages/pysss_murmur.so
usr/lib/python3.13/site-packages/pysss_nss_idmap.so
usr/lib/python3.13/site-packages/sssd/
usr/lib/python3.13/site-packages/sssd/__init__.py
usr/lib/python3.13/site-packages/sssd/modules/
usr/lib/python3.13/site-packages/sssd/modules/__init__.py
usr/lib/python3.13/site-packages/sssd/modules/error.py
usr/lib/python3.13/site-packages/sssd/modules/request.py
usr/lib/python3.13/site-packages/sssd/parser.py
usr/lib/python3.13/site-packages/sssd/source_files.py
usr/lib/python3.13/site-packages/sssd/source_journald.py
usr/lib/python3.13/site-packages/sssd/source_reader.py
usr/lib/python3.13/site-packages/sssd/sss_analyze.py
usr/lib/python3.13/site-packages/sssd/util.py
usr/lib/samba/
usr/lib/samba/idmap/
usr/lib/samba/idmap/sss.so
usr/lib/samba/ldb/
usr/lib/samba/ldb/memberof.so
usr/lib/security/
usr/lib/security/pam_sss.so
usr/lib/security/pam_sss_gss.so
usr/lib/sssd/
usr/lib/sssd/conf/
usr/lib/sssd/conf/sssd.conf
usr/lib/sssd/libifp_iface.so
usr/lib/sssd/libifp_iface_sync.so
usr/lib/sssd/libsss_ad.so
usr/lib/sssd/libsss_cert.so
usr/lib/sssd/libsss_child.so
usr/lib/sssd/libsss_crypt.so
usr/lib/sssd/libsss_debug.so
usr/lib/sssd/libsss_files.so
usr/lib/sssd/libsss_iface.so
usr/lib/sssd/libsss_iface_sync.so
usr/lib/sssd/libsss_ipa.so
usr/lib/sssd/libsss_krb5.so
usr/lib/sssd/libsss_krb5_common.so
usr/lib/sssd/libsss_ldap.so
usr/lib/sssd/libsss_ldap_common.so
usr/lib/sssd/libsss_proxy.so
usr/lib/sssd/libsss_sbus.so
usr/lib/sssd/libsss_sbus_sync.so
usr/lib/sssd/libsss_simple.so
usr/lib/sssd/libsss_util.so
usr/lib/sssd/modules/
usr/lib/sssd/modules/libsss_autofs.so
usr/lib/sssd/modules/sssd_krb5_idp_plugin.so
usr/lib/sssd/modules/sssd_krb5_localauth_plugin.so
usr/lib/sssd/sssd/
usr/lib/sssd/sssd/gpo_child
usr/lib/sssd/sssd/krb5_child
usr/lib/sssd/sssd/ldap_child
usr/lib/sssd/sssd/oidc_child
usr/lib/sssd/sssd/p11_child
usr/lib/sssd/sssd/proxy_child
usr/lib/sssd/sssd/sss_analyze
usr/lib/sssd/sssd/sss_signal
usr/lib/sssd/sssd/sssd_autofs
usr/lib/sssd/sssd/sssd_be
usr/lib/sssd/sssd/sssd_check_socket_activated_responders
usr/lib/sssd/sssd/sssd_ifp
usr/lib/sssd/sssd/sssd_kcm
usr/lib/sssd/sssd/sssd_nss
usr/lib/sssd/sssd/sssd_pac
usr/lib/sssd/sssd/sssd_pam
usr/lib/sssd/sssd/sssd_ssh
usr/lib/sssd/sssd/sssd_sudo
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sssd-autofs.service
usr/lib/systemd/system/sssd-autofs.socket
usr/lib/systemd/system/sssd-ifp.service
usr/lib/systemd/system/sssd-kcm.service
usr/lib/systemd/system/sssd-kcm.socket
usr/lib/systemd/system/sssd-nss.service
usr/lib/systemd/system/sssd-nss.socket
usr/lib/systemd/system/sssd-pac.service
usr/lib/systemd/system/sssd-pac.socket
usr/lib/systemd/system/sssd-pam.service
usr/lib/systemd/system/sssd-pam.socket
usr/lib/systemd/system/sssd-ssh.service
usr/lib/systemd/system/sssd-ssh.socket
usr/lib/systemd/system/sssd-sudo.service
usr/lib/systemd/system/sssd-sudo.socket
usr/lib/systemd/system/sssd.service
usr/lib/systemd/system/sssd.service.d/
usr/lib/systemd/system/sssd.service.d/10-perms.conf
usr/lib/sysusers.d/
usr/lib/sysusers.d/sssd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/sssd-perms.conf
usr/lib/tmpfiles.d/sssd-tmpfiles.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.sssd.infopipe.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.sssd.infopipe.conf
usr/share/locale/
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/sssd.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/sssd.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/sssd.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/sssd.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/sssd.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/sssd.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/sssd.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/sssd.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/sssd.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/sssd.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/sssd.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/sssd.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/sssd.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/sssd.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/sssd.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/sssd.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/sssd.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/sssd.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/sssd.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/sssd.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/sssd.mo
usr/share/locale/tg/
usr/share/locale/tg/LC_MESSAGES/
usr/share/locale/tg/LC_MESSAGES/sssd.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/sssd.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/sssd.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/sssd.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/sssd.mo
usr/share/man/
usr/share/man/ca/
usr/share/man/ca/man5/
usr/share/man/ca/man5/sss_rpcidmapd.5.gz
usr/share/man/ca/man5/sssd-simple.5.gz
usr/share/man/ca/man8/
usr/share/man/ca/man8/sss_obfuscate.8.gz
usr/share/man/ca/man8/sss_seed.8.gz
usr/share/man/de/
usr/share/man/de/man5/
usr/share/man/de/man5/sssd-krb5.5.gz
usr/share/man/de/man5/sssd-simple.5.gz
usr/share/man/de/man8/
usr/share/man/de/man8/sss_obfuscate.8.gz
usr/share/man/de/man8/sss_seed.8.gz
usr/share/man/es/
usr/share/man/es/man5/
usr/share/man/es/man5/sss-certmap.5.gz
usr/share/man/es/man5/sssd-ipa.5.gz
usr/share/man/es/man5/sssd-ldap-attributes.5.gz
usr/share/man/es/man5/sssd-ldap.5.gz
usr/share/man/es/man5/sssd-simple.5.gz
usr/share/man/es/man8/
usr/share/man/es/man8/sss_obfuscate.8.gz
usr/share/man/es/man8/sss_seed.8.gz
usr/share/man/es/man8/sssd_krb5_locator_plugin.8.gz
usr/share/man/fr/
usr/share/man/fr/man5/
usr/share/man/fr/man5/sssd-simple.5.gz
usr/share/man/fr/man8/
usr/share/man/fr/man8/sss_obfuscate.8.gz
usr/share/man/fr/man8/sss_seed.8.gz
usr/share/man/ja/
usr/share/man/ja/man5/
usr/share/man/ja/man5/sssd-simple.5.gz
usr/share/man/ja/man8/
usr/share/man/ja/man8/sss_obfuscate.8.gz
usr/share/man/man1/
usr/share/man/man1/sss_ssh_authorizedkeys.1.gz
usr/share/man/man1/sss_ssh_knownhosts.1.gz
usr/share/man/man5/
usr/share/man/man5/sss-certmap.5.gz
usr/share/man/man5/sss_rpcidmapd.5.gz
usr/share/man/man5/sssd-ad.5.gz
usr/share/man/man5/sssd-files.5.gz
usr/share/man/man5/sssd-ifp.5.gz
usr/share/man/man5/sssd-ipa.5.gz
usr/share/man/man5/sssd-krb5.5.gz
usr/share/man/man5/sssd-ldap-attributes.5.gz
usr/share/man/man5/sssd-ldap.5.gz
usr/share/man/man5/sssd-session-recording.5.gz
usr/share/man/man5/sssd-simple.5.gz
usr/share/man/man5/sssd-sudo.5.gz
usr/share/man/man5/sssd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/idmap_sss.8.gz
usr/share/man/man8/pam_sss.8.gz
usr/share/man/man8/pam_sss_gss.8.gz
usr/share/man/man8/sss_cache.8.gz
usr/share/man/man8/sss_debuglevel.8.gz
usr/share/man/man8/sss_obfuscate.8.gz
usr/share/man/man8/sss_override.8.gz
usr/share/man/man8/sss_seed.8.gz
usr/share/man/man8/sssctl.8.gz
usr/share/man/man8/sssd-kcm.8.gz
usr/share/man/man8/sssd.8.gz
usr/share/man/man8/sssd_krb5_localauth_plugin.8.gz
usr/share/man/man8/sssd_krb5_locator_plugin.8.gz
usr/share/man/ru/
usr/share/man/ru/man1/
usr/share/man/ru/man1/sss_ssh_authorizedkeys.1.gz
usr/share/man/ru/man1/sss_ssh_knownhostsproxy.1.gz
usr/share/man/ru/man5/
usr/share/man/ru/man5/sss-certmap.5.gz
usr/share/man/ru/man5/sss_rpcidmapd.5.gz
usr/share/man/ru/man5/sssd-ad.5.gz
usr/share/man/ru/man5/sssd-files.5.gz
usr/share/man/ru/man5/sssd-ifp.5.gz
usr/share/man/ru/man5/sssd-ipa.5.gz
usr/share/man/ru/man5/sssd-krb5.5.gz
usr/share/man/ru/man5/sssd-ldap-attributes.5.gz
usr/share/man/ru/man5/sssd-ldap.5.gz
usr/share/man/ru/man5/sssd-session-recording.5.gz
usr/share/man/ru/man5/sssd-simple.5.gz
usr/share/man/ru/man5/sssd-sudo.5.gz
usr/share/man/ru/man5/sssd-systemtap.5.gz
usr/share/man/ru/man5/sssd.conf.5.gz
usr/share/man/ru/man8/
usr/share/man/ru/man8/idmap_sss.8.gz
usr/share/man/ru/man8/pam_sss.8.gz
usr/share/man/ru/man8/pam_sss_gss.8.gz
usr/share/man/ru/man8/sss_cache.8.gz
usr/share/man/ru/man8/sss_debuglevel.8.gz
usr/share/man/ru/man8/sss_obfuscate.8.gz
usr/share/man/ru/man8/sss_override.8.gz
usr/share/man/ru/man8/sss_seed.8.gz
usr/share/man/ru/man8/sssctl.8.gz
usr/share/man/ru/man8/sssd-kcm.8.gz
usr/share/man/ru/man8/sssd.8.gz
usr/share/man/ru/man8/sssd_krb5_localauth_plugin.8.gz
usr/share/man/ru/man8/sssd_krb5_locator_plugin.8.gz
usr/share/man/sv/
usr/share/man/sv/man1/
usr/share/man/sv/man1/sss_ssh_authorizedkeys.1.gz
usr/share/man/sv/man1/sss_ssh_knownhostsproxy.1.gz
usr/share/man/sv/man5/
usr/share/man/sv/man5/sss-certmap.5.gz
usr/share/man/sv/man5/sss_rpcidmapd.5.gz
usr/share/man/sv/man5/sssd-ad.5.gz
usr/share/man/sv/man5/sssd-files.5.gz
usr/share/man/sv/man5/sssd-ifp.5.gz
usr/share/man/sv/man5/sssd-ipa.5.gz
usr/share/man/sv/man5/sssd-krb5.5.gz
usr/share/man/sv/man5/sssd-ldap-attributes.5.gz
usr/share/man/sv/man5/sssd-ldap.5.gz
usr/share/man/sv/man5/sssd-session-recording.5.gz
usr/share/man/sv/man5/sssd-simple.5.gz
usr/share/man/sv/man5/sssd-sudo.5.gz
usr/share/man/sv/man5/sssd-systemtap.5.gz
usr/share/man/sv/man5/sssd.conf.5.gz
usr/share/man/sv/man8/
usr/share/man/sv/man8/idmap_sss.8.gz
usr/share/man/sv/man8/pam_sss.8.gz
usr/share/man/sv/man8/pam_sss_gss.8.gz
usr/share/man/sv/man8/sss_cache.8.gz
usr/share/man/sv/man8/sss_debuglevel.8.gz
usr/share/man/sv/man8/sss_obfuscate.8.gz
usr/share/man/sv/man8/sss_override.8.gz
usr/share/man/sv/man8/sss_seed.8.gz
usr/share/man/sv/man8/sssctl.8.gz
usr/share/man/sv/man8/sssd-kcm.8.gz
usr/share/man/sv/man8/sssd.8.gz
usr/share/man/sv/man8/sssd_krb5_localauth_plugin.8.gz
usr/share/man/sv/man8/sssd_krb5_locator_plugin.8.gz
usr/share/man/uk/
usr/share/man/uk/man1/
usr/share/man/uk/man1/sss_ssh_authorizedkeys.1.gz
usr/share/man/uk/man1/sss_ssh_knownhostsproxy.1.gz
usr/share/man/uk/man5/
usr/share/man/uk/man5/sss-certmap.5.gz
usr/share/man/uk/man5/sss_rpcidmapd.5.gz
usr/share/man/uk/man5/sssd-ad.5.gz
usr/share/man/uk/man5/sssd-files.5.gz
usr/share/man/uk/man5/sssd-ifp.5.gz
usr/share/man/uk/man5/sssd-ipa.5.gz
usr/share/man/uk/man5/sssd-krb5.5.gz
usr/share/man/uk/man5/sssd-ldap-attributes.5.gz
usr/share/man/uk/man5/sssd-ldap.5.gz
usr/share/man/uk/man5/sssd-session-recording.5.gz
usr/share/man/uk/man5/sssd-simple.5.gz
usr/share/man/uk/man5/sssd-sudo.5.gz
usr/share/man/uk/man5/sssd-systemtap.5.gz
usr/share/man/uk/man5/sssd.conf.5.gz
usr/share/man/uk/man8/
usr/share/man/uk/man8/idmap_sss.8.gz
usr/share/man/uk/man8/pam_sss.8.gz
usr/share/man/uk/man8/pam_sss_gss.8.gz
usr/share/man/uk/man8/sss_cache.8.gz
usr/share/man/uk/man8/sss_debuglevel.8.gz
usr/share/man/uk/man8/sss_obfuscate.8.gz
usr/share/man/uk/man8/sss_override.8.gz
usr/share/man/uk/man8/sss_seed.8.gz
usr/share/man/uk/man8/sssctl.8.gz
usr/share/man/uk/man8/sssd-kcm.8.gz
usr/share/man/uk/man8/sssd.8.gz
usr/share/man/uk/man8/sssd_krb5_localauth_plugin.8.gz
usr/share/man/uk/man8/sssd_krb5_locator_plugin.8.gz
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/sssd-pcsc.rules
usr/share/sssd-kcm/
usr/share/sssd-kcm/kcm_default_ccache
usr/share/sssd/
usr/share/sssd/cfg_rules.ini
usr/share/sssd/krb5-snippets/
usr/share/sssd/krb5-snippets/enable_sssd_conf_dir
usr/share/sssd/krb5-snippets/sssd_enable_idp
usr/share/sssd/sssd.api.conf
usr/share/sssd/sssd.api.d/
usr/share/sssd/sssd.api.d/sssd-ad.conf
usr/share/sssd/sssd.api.d/sssd-files.conf
usr/share/sssd/sssd.api.d/sssd-ipa.conf
usr/share/sssd/sssd.api.d/sssd-krb5.conf
usr/share/sssd/sssd.api.d/sssd-ldap.conf
usr/share/sssd/sssd.api.d/sssd-proxy.conf
usr/share/sssd/sssd.api.d/sssd-simple.conf
var/
var/lib/
var/lib/sss/
var/lib/sss/db/
var/lib/sss/deskprofile/
var/lib/sss/gpo_cache/
var/lib/sss/keytabs/
var/lib/sss/mc/
var/lib/sss/pipes/
var/lib/sss/pipes/private/
var/lib/sss/pubconf/
var/lib/sss/pubconf/krb5.include.d/
var/lib/sss/secrets/
var/log/
var/log/sssd/
