# extracted/sslh-2.1.2-1/files
%FILES%
etc/
etc/sslh.cfg
usr/
usr/bin/
usr/bin/sslh
usr/bin/sslh-ev
usr/bin/sslh-fork
usr/bin/sslh-select
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/systemd-sslh-generator
usr/lib/systemd/system/
usr/lib/systemd/system/sslh-fork.service
usr/lib/systemd/system/sslh-select.service
usr/lib/systemd/system/sslh.service
usr/share/
usr/share/doc/
usr/share/doc/sslh/
usr/share/doc/sslh/basic.cfg
usr/share/doc/sslh/example.cfg
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/sslh-ev.8.gz
usr/share/man/man8/sslh-fork.8.gz
usr/share/man/man8/sslh-select.8.gz
usr/share/man/man8/sslh.8.gz
