# extracted/mimir-2.12.0-1/files
%FILES%
etc/
etc/mimir.yaml
usr/
usr/bin/
usr/bin/metaconvert
usr/bin/mimir
usr/bin/mimir-continuous-test
usr/bin/mimirtool
usr/bin/query-tee
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mimir.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/mimir.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mimir.conf
