# extracted/prometheus-smartctl-exporter-0.12.0-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-smartctl-exporter
usr/
usr/bin/
usr/bin/prometheus-smartctl-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-smartctl-exporter.service
usr/share/
usr/share/licenses/
usr/share/licenses/prometheus-smartctl-exporter/
usr/share/licenses/prometheus-smartctl-exporter/LICENSE
