# extracted/meshbird-2.3-7/files
%FILES%
etc/
etc/meshbird.conf
usr/
usr/bin/
usr/bin/meshbird
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/meshbird.service
