# extracted/pdnsd-1.2.9.a-15/files
%FILES%
usr/
usr/bin/
usr/bin/pdnsd
usr/bin/pdnsd-ctl
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pdnsd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pdnsd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pdnsd.conf
usr/share/
usr/share/doc/
usr/share/doc/pdnsd/
usr/share/doc/pdnsd/pdnsd.conf
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/pdnsd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/pdnsd-ctl.8.gz
usr/share/man/man8/pdnsd.8.gz
