# extracted/deepin-anything-5.0.18-1/files
%FILES%
etc/
etc/dbus-1/
usr/
usr/bin/
usr/bin/deepin-anything-monitor
usr/bin/deepin-anything-tool
usr/bin/deepin-anything-tool-ionice
usr/include/
usr/include/deepin-anything-server/
usr/include/deepin-anything-server/dasdefine.h
usr/include/deepin-anything-server/dasfactory.h
usr/include/deepin-anything-server/dasinterface.h
usr/include/deepin-anything-server/dasplugin.h
usr/include/deepin-anything-server/lftmanager.h
usr/include/deepin-anything/
usr/include/deepin-anything/fs_buf.h
usr/include/deepin-anything/index/
usr/include/deepin-anything/index/composite_str.h
usr/include/deepin-anything/index/index.h
usr/include/deepin-anything/index/index_allfile.h
usr/include/deepin-anything/index/index_allmem.h
usr/include/deepin-anything/index/index_base.h
usr/include/deepin-anything/index/index_utils.h
usr/include/deepin-anything/utils.h
usr/include/deepin-anything/vfs_change_consts.h
usr/include/deepin-anything/vfs_change_uapi.h
usr/include/deepin-anything/walkdir.h
usr/lib/
usr/lib/deepin-anything-server-lib/
usr/lib/deepin-anything-server-lib/plugins/
usr/lib/deepin-anything-server-lib/plugins/handlers/
usr/lib/deepin-anything-server-lib/plugins/handlers/README.txt
usr/lib/deepin-anything-server-lib/plugins/handlers/libupdate-lft.so
usr/lib/libanything.so
usr/lib/libanything.so.1
usr/lib/libanything.so.1.0
usr/lib/libanything.so.1.0.0
usr/lib/libdeepin-anything-server-lib.so
usr/lib/libdeepin-anything-server-lib.so.0
usr/lib/libdeepin-anything-server-lib.so.0.0
usr/lib/libdeepin-anything-server-lib.so.0.0.1
usr/lib/modules-load.d/
usr/lib/modules-load.d/anything.conf
usr/lib/pkgconfig/
usr/lib/pkgconfig/deepin-anything-server-lib.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/deepin-anything-monitor.service
usr/lib/systemd/system/deepin-anything-tool.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/deepin-anything-monitor.conf
usr/lib/sysusers.d/deepin-anything-server.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/com.deepin.anything.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/com.deepin.anything.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.deepin.anything.conf
