# extracted/mpd-0.23.16-1/files
%FILES%
etc/
etc/mpd.conf
usr/
usr/bin/
usr/bin/mpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/mpd.service
usr/lib/systemd/system/mpd.service.d/
usr/lib/systemd/system/mpd.service.d/00-arch.conf
usr/lib/systemd/system/mpd.socket
usr/lib/systemd/user/
usr/lib/systemd/user/mpd.service
usr/lib/systemd/user/mpd.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/mpd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/mpd.conf
usr/share/
usr/share/doc/
usr/share/doc/mpd/
usr/share/doc/mpd/AUTHORS
usr/share/doc/mpd/COPYING
usr/share/doc/mpd/NEWS
usr/share/doc/mpd/README.md
usr/share/doc/mpd/html/
usr/share/doc/mpd/html/.buildinfo
usr/share/doc/mpd/html/_sources/
usr/share/doc/mpd/html/_sources/client.rst.txt
usr/share/doc/mpd/html/_sources/developer.rst.txt
usr/share/doc/mpd/html/_sources/index.rst.txt
usr/share/doc/mpd/html/_sources/mpd.1.rst.txt
usr/share/doc/mpd/html/_sources/mpd.conf.5.rst.txt
usr/share/doc/mpd/html/_sources/plugins.rst.txt
usr/share/doc/mpd/html/_sources/protocol.rst.txt
usr/share/doc/mpd/html/_sources/user.rst.txt
usr/share/doc/mpd/html/_static/
usr/share/doc/mpd/html/_static/basic.css
usr/share/doc/mpd/html/_static/classic.css
usr/share/doc/mpd/html/_static/doctools.js
usr/share/doc/mpd/html/_static/documentation_options.js
usr/share/doc/mpd/html/_static/file.png
usr/share/doc/mpd/html/_static/language_data.js
usr/share/doc/mpd/html/_static/minus.png
usr/share/doc/mpd/html/_static/plus.png
usr/share/doc/mpd/html/_static/pygments.css
usr/share/doc/mpd/html/_static/searchtools.js
usr/share/doc/mpd/html/_static/sidebar.js
usr/share/doc/mpd/html/_static/sphinx_highlight.js
usr/share/doc/mpd/html/client.html
usr/share/doc/mpd/html/developer.html
usr/share/doc/mpd/html/genindex.html
usr/share/doc/mpd/html/index.html
usr/share/doc/mpd/html/mpd.1.html
usr/share/doc/mpd/html/mpd.conf.5.html
usr/share/doc/mpd/html/objects.inv
usr/share/doc/mpd/html/plugins.html
usr/share/doc/mpd/html/protocol.html
usr/share/doc/mpd/html/search.html
usr/share/doc/mpd/html/searchindex.js
usr/share/doc/mpd/html/user.html
usr/share/doc/mpd/mpdconf.example
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/mpd.svg
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mpd.1.gz
usr/share/man/man5/
usr/share/man/man5/mpd.conf.5.gz
