# extracted/miniupnpd-2.3.7-1/files
%FILES%
etc/
etc/miniupnpd/
etc/miniupnpd/ip6tables_init.sh
etc/miniupnpd/ip6tables_removeall.sh
etc/miniupnpd/iptables_init.sh
etc/miniupnpd/iptables_removeall.sh
etc/miniupnpd/miniupnpd.conf
etc/miniupnpd/miniupnpd_functions.sh
usr/
usr/bin/
usr/bin/miniupnpd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/miniupnpd.service
usr/share/
usr/share/licenses/
usr/share/licenses/miniupnpd/
usr/share/licenses/miniupnpd/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/miniupnpd.8.gz
