# extracted/nginx-1.26.2-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/nginx
etc/nginx/
etc/nginx/fastcgi.conf
etc/nginx/fastcgi_params
etc/nginx/koi-utf
etc/nginx/koi-win
etc/nginx/nginx.conf
etc/nginx/scgi_params
etc/nginx/uwsgi_params
etc/nginx/win-utf
usr/
usr/bin/
usr/bin/nginx
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/nginx.service
usr/share/
usr/share/licenses/
usr/share/licenses/nginx/
usr/share/licenses/nginx/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/nginx.8.gz
usr/share/nginx/
usr/share/nginx/html/
usr/share/nginx/html/50x.html
usr/share/nginx/html/index.html
usr/share/vim/
usr/share/vim/vimfiles/
usr/share/vim/vimfiles/ftdetect/
usr/share/vim/vimfiles/ftdetect/nginx.vim
usr/share/vim/vimfiles/ftplugin/
usr/share/vim/vimfiles/ftplugin/nginx.vim
usr/share/vim/vimfiles/indent/
usr/share/vim/vimfiles/indent/nginx.vim
usr/share/vim/vimfiles/syntax/
usr/share/vim/vimfiles/syntax/nginx.vim
var/
var/lib/
var/lib/nginx/
var/lib/nginx/proxy/
var/log/
var/log/nginx/
