# extracted/lynis-3.1.2-1/files
%FILES%
etc/
etc/lynis/
etc/lynis/default.prf
usr/
usr/bin/
usr/bin/lynis
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/lynis.service
usr/lib/systemd/system/lynis.timer
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/lynis
usr/share/doc/
usr/share/doc/lynis/
usr/share/doc/lynis/CHANGELOG.md
usr/share/doc/lynis/FAQ
usr/share/doc/lynis/INSTALL
usr/share/doc/lynis/README
usr/share/lynis/
usr/share/lynis/db/
usr/share/lynis/db/fileperms.db
usr/share/lynis/db/hints.db
usr/share/lynis/db/integrity.db
usr/share/lynis/db/languages/
usr/share/lynis/db/languages/az
usr/share/lynis/db/languages/br
usr/share/lynis/db/languages/cn
usr/share/lynis/db/languages/da
usr/share/lynis/db/languages/de
usr/share/lynis/db/languages/de-AT
usr/share/lynis/db/languages/en
usr/share/lynis/db/languages/en-GB
usr/share/lynis/db/languages/en-US
usr/share/lynis/db/languages/es
usr/share/lynis/db/languages/fi
usr/share/lynis/db/languages/fr
usr/share/lynis/db/languages/gr
usr/share/lynis/db/languages/he
usr/share/lynis/db/languages/hu
usr/share/lynis/db/languages/id
usr/share/lynis/db/languages/it
usr/share/lynis/db/languages/ja
usr/share/lynis/db/languages/ko
usr/share/lynis/db/languages/nb-NO
usr/share/lynis/db/languages/nl
usr/share/lynis/db/languages/nl-BE
usr/share/lynis/db/languages/nl-NL
usr/share/lynis/db/languages/pl
usr/share/lynis/db/languages/pt
usr/share/lynis/db/languages/ru
usr/share/lynis/db/languages/se
usr/share/lynis/db/languages/sk
usr/share/lynis/db/languages/tr
usr/share/lynis/db/malware-susp.db
usr/share/lynis/db/malware.db
usr/share/lynis/db/sbl.db
usr/share/lynis/db/software-eol.db
usr/share/lynis/db/tests.db
usr/share/lynis/include/
usr/share/lynis/include/binaries
usr/share/lynis/include/consts
usr/share/lynis/include/data_upload
usr/share/lynis/include/functions
usr/share/lynis/include/helper_audit_dockerfile
usr/share/lynis/include/helper_configure
usr/share/lynis/include/helper_generate
usr/share/lynis/include/helper_show
usr/share/lynis/include/helper_system_remote_scan
usr/share/lynis/include/helper_update
usr/share/lynis/include/osdetection
usr/share/lynis/include/parameters
usr/share/lynis/include/profiles
usr/share/lynis/include/report
usr/share/lynis/include/tests_accounting
usr/share/lynis/include/tests_authentication
usr/share/lynis/include/tests_banners
usr/share/lynis/include/tests_boot_services
usr/share/lynis/include/tests_containers
usr/share/lynis/include/tests_crypto
usr/share/lynis/include/tests_custom.template
usr/share/lynis/include/tests_databases
usr/share/lynis/include/tests_dns
usr/share/lynis/include/tests_file_integrity
usr/share/lynis/include/tests_file_permissions
usr/share/lynis/include/tests_filesystems
usr/share/lynis/include/tests_firewalls
usr/share/lynis/include/tests_hardening
usr/share/lynis/include/tests_homedirs
usr/share/lynis/include/tests_insecure_services
usr/share/lynis/include/tests_kerberos
usr/share/lynis/include/tests_kernel
usr/share/lynis/include/tests_kernel_hardening
usr/share/lynis/include/tests_ldap
usr/share/lynis/include/tests_logging
usr/share/lynis/include/tests_mac_frameworks
usr/share/lynis/include/tests_mail_messaging
usr/share/lynis/include/tests_malware
usr/share/lynis/include/tests_memory_processes
usr/share/lynis/include/tests_nameservices
usr/share/lynis/include/tests_networking
usr/share/lynis/include/tests_php
usr/share/lynis/include/tests_ports_packages
usr/share/lynis/include/tests_printers_spoolers
usr/share/lynis/include/tests_scheduling
usr/share/lynis/include/tests_shells
usr/share/lynis/include/tests_snmp
usr/share/lynis/include/tests_squid
usr/share/lynis/include/tests_ssh
usr/share/lynis/include/tests_storage
usr/share/lynis/include/tests_storage_nfs
usr/share/lynis/include/tests_system_integrity
usr/share/lynis/include/tests_time
usr/share/lynis/include/tests_tooling
usr/share/lynis/include/tests_usb
usr/share/lynis/include/tests_virtualization
usr/share/lynis/include/tests_webservers
usr/share/lynis/include/tool_tips
usr/share/lynis/plugins/
usr/share/lynis/plugins/README
usr/share/lynis/plugins/custom_plugin.template
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/lynis.8.gz
