# extracted/prometheus-snmp-exporter-0.26.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-snmp-exporter
usr/
usr/bin/
usr/bin/prometheus-snmp-exporter
usr/bin/prometheus-snmp-generator
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-snmp-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-snmp-exporter.conf
usr/share/
usr/share/doc/
usr/share/doc/prometheus-snmp-exporter/
usr/share/doc/prometheus-snmp-exporter/README.md
usr/share/doc/prometheus-snmp-exporter/auth-split-migration.md
usr/share/doc/prometheus-snmp-exporter/generator/
usr/share/doc/prometheus-snmp-exporter/generator/README.md
usr/share/doc/prometheus-snmp-exporter/generator/generator.yml
usr/share/doc/prometheus-snmp-exporter/snmp_example.yml
usr/share/licenses/
usr/share/licenses/prometheus-snmp-exporter/
usr/share/licenses/prometheus-snmp-exporter/LICENSE
