# extracted/tftp-hpa-5.2-10/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/tftpd
srv/
srv/tftp/
usr/
usr/bin/
usr/bin/in.tftpd
usr/bin/tftp
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tftpd.service
usr/lib/systemd/system/tftpd.socket
usr/share/
usr/share/licenses/
usr/share/licenses/tftp-hpa/
usr/share/licenses/tftp-hpa/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/tftp.1.gz
usr/share/man/man8/
usr/share/man/man8/in.tftpd.8.gz
