# extracted/oidentd-3.1.0-1/files
%FILES%
etc/
etc/oidentd.conf
etc/oidentd_masq.conf
usr/
usr/bin/
usr/bin/oidentd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/oidentd.socket
usr/lib/systemd/system/oidentd@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/oidentd.conf
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/oidentd.conf.5.gz
usr/share/man/man5/oidentd_masq.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/oidentd.8.gz
