# extracted/arch-release-promotion-0.3.0-2/files
%FILES%
etc/
etc/arch-release-promotion/
etc/arch-release-promotion/projects.toml
usr/
usr/bin/
usr/bin/arch-release-promotion
usr/bin/arch-release-sync
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/METADATA
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/RECORD
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/WHEEL
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/entry_points.txt
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/licenses/
usr/lib/python3.12/site-packages/arch_release_promotion-0.3.0.dist-info/licenses/LICENSE
usr/lib/python3.12/site-packages/arch_release_promotion/
usr/lib/python3.12/site-packages/arch_release_promotion/__init__.py
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/__init__.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/__init__.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/argparse.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/argparse.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/cli.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/cli.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/config.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/config.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/files.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/files.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/gitlab.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/gitlab.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/release.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/release.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/signature.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/signature.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/torrent.cpython-312.opt-1.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/__pycache__/torrent.cpython-312.pyc
usr/lib/python3.12/site-packages/arch_release_promotion/argparse.py
usr/lib/python3.12/site-packages/arch_release_promotion/cli.py
usr/lib/python3.12/site-packages/arch_release_promotion/config.py
usr/lib/python3.12/site-packages/arch_release_promotion/files.py
usr/lib/python3.12/site-packages/arch_release_promotion/gitlab.py
usr/lib/python3.12/site-packages/arch_release_promotion/release.py
usr/lib/python3.12/site-packages/arch_release_promotion/signature.py
usr/lib/python3.12/site-packages/arch_release_promotion/torrent.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/arch-release-sync.service
usr/lib/systemd/system/arch-release-sync.timer
usr/lib/sysusers.d/
usr/lib/sysusers.d/arch-release-promotion.conf
usr/share/
usr/share/doc/
usr/share/doc/arch-release-promotion/
usr/share/doc/arch-release-promotion/README.rst
usr/share/doc/arch-release-promotion/example.toml
