# extracted/hdapsd-20141203-5/files
%FILES%
usr/
usr/bin/
usr/bin/hdapsd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hdapsd@.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/hdapsd.rules
usr/share/
usr/share/doc/
usr/share/doc/hdapsd/
usr/share/doc/hdapsd/ChangeLog
usr/share/doc/hdapsd/README.md
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/hdapsd.8.gz
