# extracted/qemu-pr-helper-9.2.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/qemu-pr-helper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/qemu-pr-helper.service
usr/lib/systemd/system/qemu-pr-helper.socket
usr/share/
usr/share/licenses/
usr/share/licenses/qemu-pr-helper/
usr/share/licenses/qemu-pr-helper/BSD-2-Clause.txt
usr/share/licenses/qemu-pr-helper/LICENSE
usr/share/licenses/qemu-pr-helper/audio.LICENSE.ISC.txt
usr/share/licenses/qemu-pr-helper/edk2-licenses.txt
usr/share/licenses/qemu-pr-helper/fuse.LICENSE.BSD-2-Clause.txt
usr/share/licenses/qemu-pr-helper/tcg-arm.LICENSE.MIT.txt
usr/share/licenses/qemu-pr-helper/tcg.LICENSE.MIT.txt
usr/share/licenses/qemu-pr-helper/tci.LICENSE.MIT.txt
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/qemu-pr-helper.8.gz
