# extracted/soju-0.8.2-2/files
%FILES%
etc/
etc/soju/
etc/soju/config
usr/
usr/bin/
usr/bin/soju
usr/bin/soju-casemap-logs
usr/bin/soju-migrate-db
usr/bin/soju-migrate-logs
usr/bin/soju-znc-import
usr/bin/sojuctl
usr/bin/sojudb
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/soju.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/soju.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/soju.conf
usr/share/
usr/share/doc/
usr/share/doc/soju/
usr/share/doc/soju/README.md
usr/share/doc/soju/clients.md
usr/share/doc/soju/doc/
usr/share/doc/soju/doc/architecture.md
usr/share/doc/soju/doc/ext/
usr/share/doc/soju/doc/ext/README.md
usr/share/doc/soju/doc/ext/account-required.md
usr/share/doc/soju/doc/ext/bouncer-networks.md
usr/share/doc/soju/doc/ext/filehost.md
usr/share/doc/soju/doc/ext/no-implicit-names.md
usr/share/doc/soju/doc/ext/read.md
usr/share/doc/soju/doc/ext/search.md
usr/share/doc/soju/doc/ext/webpush.md
usr/share/doc/soju/doc/getting-started.md
usr/share/doc/soju/doc/packaging.md
usr/share/doc/soju/doc/per-user-ip.md
usr/share/doc/soju/doc/sojuctl.1.scd
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/soju.1.gz
usr/share/man/man1/sojuctl.1.gz
