# extracted/dehydrated-0.7.1-2/files
%FILES%
etc/
etc/dehydrated/
usr/
usr/bin/
usr/bin/dehydrated
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dehydrated.service
usr/lib/systemd/system/dehydrated.timer
usr/share/
usr/share/doc/
usr/share/doc/dehydrated/
usr/share/doc/dehydrated/acme-v1.md
usr/share/doc/dehydrated/dns-verification.md
usr/share/doc/dehydrated/domains_txt.md
usr/share/doc/dehydrated/ecc.md
usr/share/doc/dehydrated/examples/
usr/share/doc/dehydrated/examples/config
usr/share/doc/dehydrated/examples/domains.txt
usr/share/doc/dehydrated/examples/hook.sh
usr/share/doc/dehydrated/hook_chain.md
usr/share/doc/dehydrated/per-certificate-config.md
usr/share/doc/dehydrated/staging.md
usr/share/doc/dehydrated/tls-alpn.md
usr/share/doc/dehydrated/troubleshooting.md
usr/share/doc/dehydrated/wellknown.md
usr/share/licenses/
usr/share/licenses/dehydrated/
usr/share/licenses/dehydrated/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dehydrated.1.gz
