# extracted/hyperv-6.12-8/files
%FILES%
usr/
usr/bin/
usr/bin/hv_fcopy_uio_daemon
usr/bin/hv_kvp_daemon
usr/bin/hv_vss_daemon
usr/bin/lsvmbus
usr/lib/
usr/lib/hypervkvpd/
usr/lib/hypervkvpd/hv_get_dhcp_info
usr/lib/hypervkvpd/hv_get_dns_info
usr/lib/hypervkvpd/hv_set_ifconfig
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hv_kvp_daemon.service
usr/lib/systemd/system/hv_vss_daemon.service
var/
var/lib/
