# extracted/openvswitch-3.3.0-1/files
%FILES%
etc/
etc/openvswitch/
usr/
usr/bin/
usr/bin/ovs-appctl
usr/bin/ovs-bugtool
usr/bin/ovs-docker
usr/bin/ovs-dpctl
usr/bin/ovs-dpctl-top
usr/bin/ovs-l3ping
usr/bin/ovs-ofctl
usr/bin/ovs-parse-backtrace
usr/bin/ovs-pcap
usr/bin/ovs-pki
usr/bin/ovs-tcpdump
usr/bin/ovs-tcpundump
usr/bin/ovs-test
usr/bin/ovs-testcontroller
usr/bin/ovs-vlan-test
usr/bin/ovs-vsctl
usr/bin/ovs-vswitchd
usr/bin/ovsdb-client
usr/bin/ovsdb-server
usr/bin/ovsdb-tool
usr/bin/vtep-ctl
usr/include/
usr/include/openflow/
usr/include/openflow/intel-ext.h
usr/include/openflow/netronome-ext.h
usr/include/openflow/nicira-ext.h
usr/include/openflow/openflow-1.0.h
usr/include/openflow/openflow-1.1.h
usr/include/openflow/openflow-1.2.h
usr/include/openflow/openflow-1.3.h
usr/include/openflow/openflow-1.4.h
usr/include/openflow/openflow-1.5.h
usr/include/openflow/openflow-common.h
usr/include/openflow/openflow.h
usr/include/openvswitch/
usr/include/openvswitch/compiler.h
usr/include/openvswitch/dynamic-string.h
usr/include/openvswitch/flow.h
usr/include/openvswitch/geneve.h
usr/include/openvswitch/hmap.h
usr/include/openvswitch/json.h
usr/include/openvswitch/list.h
usr/include/openvswitch/match.h
usr/include/openvswitch/meta-flow.h
usr/include/openvswitch/namemap.h
usr/include/openvswitch/netdev.h
usr/include/openvswitch/nsh.h
usr/include/openvswitch/ofp-actions.h
usr/include/openvswitch/ofp-bundle.h
usr/include/openvswitch/ofp-connection.h
usr/include/openvswitch/ofp-ct.h
usr/include/openvswitch/ofp-ed-props.h
usr/include/openvswitch/ofp-errors.h
usr/include/openvswitch/ofp-flow.h
usr/include/openvswitch/ofp-group.h
usr/include/openvswitch/ofp-ipfix.h
usr/include/openvswitch/ofp-match.h
usr/include/openvswitch/ofp-meter.h
usr/include/openvswitch/ofp-monitor.h
usr/include/openvswitch/ofp-msgs.h
usr/include/openvswitch/ofp-packet.h
usr/include/openvswitch/ofp-parse.h
usr/include/openvswitch/ofp-port.h
usr/include/openvswitch/ofp-print.h
usr/include/openvswitch/ofp-prop.h
usr/include/openvswitch/ofp-protocol.h
usr/include/openvswitch/ofp-queue.h
usr/include/openvswitch/ofp-switch.h
usr/include/openvswitch/ofp-table.h
usr/include/openvswitch/ofp-util.h
usr/include/openvswitch/ofpbuf.h
usr/include/openvswitch/packets.h
usr/include/openvswitch/poll-loop.h
usr/include/openvswitch/rconn.h
usr/include/openvswitch/shash.h
usr/include/openvswitch/thread.h
usr/include/openvswitch/token-bucket.h
usr/include/openvswitch/tun-metadata.h
usr/include/openvswitch/type-props.h
usr/include/openvswitch/types.h
usr/include/openvswitch/usdt-probes.h
usr/include/openvswitch/util.h
usr/include/openvswitch/uuid.h
usr/include/openvswitch/vconn.h
usr/include/openvswitch/version.h
usr/include/openvswitch/vlog.h
usr/lib/
usr/lib/libofproto.a
usr/lib/libopenvswitch.a
usr/lib/libopenvswitchavx512.a
usr/lib/libovsdb.a
usr/lib/libsflow.a
usr/lib/libvtep.a
usr/lib/pkgconfig/
usr/lib/pkgconfig/libofproto.pc
usr/lib/pkgconfig/libopenvswitch.pc
usr/lib/pkgconfig/libovsdb.pc
usr/lib/pkgconfig/libsflow.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ovs-vswitchd.service
usr/lib/systemd/system/ovsdb-server.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/openvswitch.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/ovs-appctl-bashcomp.bash
usr/share/bash-completion/completions/ovs-vsctl-bashcomp.bash
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/ovs-pcap.1.gz
usr/share/man/man1/ovsdb-client.1.gz
usr/share/man/man1/ovsdb-server.1.gz
usr/share/man/man1/ovsdb-tool.1.gz
usr/share/man/man5/
usr/share/man/man5/ovs-vswitchd.conf.db.5.gz
usr/share/man/man5/ovsdb-server.5.gz
usr/share/man/man5/ovsdb.local-config.5.gz
usr/share/man/man5/vtep.5.gz
usr/share/man/man7/
usr/share/man/man7/ovs-fields.7.gz
usr/share/man/man8/
usr/share/man/man8/ovs-bugtool.8.gz
usr/share/man/man8/ovs-dpctl-top.8.gz
usr/share/man/man8/ovs-dpctl.8.gz
usr/share/man/man8/ovs-kmod-ctl.8.gz
usr/share/man/man8/ovs-ofctl.8.gz
usr/share/man/man8/ovs-testcontroller.8.gz
usr/share/man/man8/ovs-vsctl.8.gz
usr/share/man/man8/ovs-vswitchd.8.gz
usr/share/man/man8/vtep-ctl.8.gz
usr/share/openvswitch/
usr/share/openvswitch/bugtool-plugins/
usr/share/openvswitch/bugtool-plugins/kernel-info/
usr/share/openvswitch/bugtool-plugins/kernel-info/openvswitch.xml
usr/share/openvswitch/bugtool-plugins/network-status/
usr/share/openvswitch/bugtool-plugins/network-status/openvswitch.xml
usr/share/openvswitch/bugtool-plugins/system-configuration.xml
usr/share/openvswitch/bugtool-plugins/system-configuration/
usr/share/openvswitch/bugtool-plugins/system-configuration/openvswitch.xml
usr/share/openvswitch/bugtool-plugins/system-logs/
usr/share/openvswitch/bugtool-plugins/system-logs/openvswitch.xml
usr/share/openvswitch/local-config.ovsschema
usr/share/openvswitch/python/
usr/share/openvswitch/python/ovs/
usr/share/openvswitch/python/ovs/__init__.py
usr/share/openvswitch/python/ovs/compat/
usr/share/openvswitch/python/ovs/compat/__init__.py
usr/share/openvswitch/python/ovs/compat/sortedcontainers/
usr/share/openvswitch/python/ovs/compat/sortedcontainers/__init__.py
usr/share/openvswitch/python/ovs/compat/sortedcontainers/sorteddict.py
usr/share/openvswitch/python/ovs/compat/sortedcontainers/sortedlist.py
usr/share/openvswitch/python/ovs/compat/sortedcontainers/sortedset.py
usr/share/openvswitch/python/ovs/daemon.py
usr/share/openvswitch/python/ovs/db/
usr/share/openvswitch/python/ovs/db/__init__.py
usr/share/openvswitch/python/ovs/db/custom_index.py
usr/share/openvswitch/python/ovs/db/data.py
usr/share/openvswitch/python/ovs/db/error.py
usr/share/openvswitch/python/ovs/db/idl.py
usr/share/openvswitch/python/ovs/db/parser.py
usr/share/openvswitch/python/ovs/db/schema.py
usr/share/openvswitch/python/ovs/db/types.py
usr/share/openvswitch/python/ovs/dirs.py
usr/share/openvswitch/python/ovs/dns_resolve.py
usr/share/openvswitch/python/ovs/fatal_signal.py
usr/share/openvswitch/python/ovs/fcntl_win.py
usr/share/openvswitch/python/ovs/flow/
usr/share/openvswitch/python/ovs/flow/__init__.py
usr/share/openvswitch/python/ovs/flow/decoders.py
usr/share/openvswitch/python/ovs/flow/filter.py
usr/share/openvswitch/python/ovs/flow/flow.py
usr/share/openvswitch/python/ovs/flow/kv.py
usr/share/openvswitch/python/ovs/flow/list.py
usr/share/openvswitch/python/ovs/flow/odp.py
usr/share/openvswitch/python/ovs/flow/ofp.py
usr/share/openvswitch/python/ovs/flow/ofp_act.py
usr/share/openvswitch/python/ovs/flow/ofp_fields.py
usr/share/openvswitch/python/ovs/json.py
usr/share/openvswitch/python/ovs/jsonrpc.py
usr/share/openvswitch/python/ovs/ovsuuid.py
usr/share/openvswitch/python/ovs/poller.py
usr/share/openvswitch/python/ovs/process.py
usr/share/openvswitch/python/ovs/reconnect.py
usr/share/openvswitch/python/ovs/socket_util.py
usr/share/openvswitch/python/ovs/stream.py
usr/share/openvswitch/python/ovs/timeval.py
usr/share/openvswitch/python/ovs/unixctl/
usr/share/openvswitch/python/ovs/unixctl/__init__.py
usr/share/openvswitch/python/ovs/unixctl/client.py
usr/share/openvswitch/python/ovs/unixctl/server.py
usr/share/openvswitch/python/ovs/util.py
usr/share/openvswitch/python/ovs/version.py
usr/share/openvswitch/python/ovs/vlog.py
usr/share/openvswitch/python/ovs/winutils.py
usr/share/openvswitch/python/ovstest/
usr/share/openvswitch/python/ovstest/__init__.py
usr/share/openvswitch/python/ovstest/args.py
usr/share/openvswitch/python/ovstest/rpcserver.py
usr/share/openvswitch/python/ovstest/tcp.py
usr/share/openvswitch/python/ovstest/tests.py
usr/share/openvswitch/python/ovstest/udp.py
usr/share/openvswitch/python/ovstest/util.py
usr/share/openvswitch/python/ovstest/vswitch.py
usr/share/openvswitch/scripts/
usr/share/openvswitch/scripts/ovs-bugtool-daemons-ver
usr/share/openvswitch/scripts/ovs-bugtool-fdb-show
usr/share/openvswitch/scripts/ovs-bugtool-get-dpdk-nic-numa
usr/share/openvswitch/scripts/ovs-bugtool-get-port-stats
usr/share/openvswitch/scripts/ovs-bugtool-ovs-appctl-dpif
usr/share/openvswitch/scripts/ovs-bugtool-ovs-bridge-datapath-type
usr/share/openvswitch/scripts/ovs-bugtool-ovs-ofctl-loop-over-bridges
usr/share/openvswitch/scripts/ovs-bugtool-ovs-vswitchd-threads-affinity
usr/share/openvswitch/scripts/ovs-bugtool-qos-configs
usr/share/openvswitch/scripts/ovs-bugtool-tc-class-show
usr/share/openvswitch/scripts/ovs-check-dead-ifs
usr/share/openvswitch/scripts/ovs-ctl
usr/share/openvswitch/scripts/ovs-kmod-ctl
usr/share/openvswitch/scripts/ovs-lib
usr/share/openvswitch/scripts/ovs-monitor-ipsec
usr/share/openvswitch/scripts/ovs-save
usr/share/openvswitch/scripts/ovs-vtep
usr/share/openvswitch/scripts/usdt/
usr/share/openvswitch/scripts/usdt/bridge_loop.bt
usr/share/openvswitch/scripts/usdt/dpif_nl_exec_monitor.py
usr/share/openvswitch/scripts/usdt/kernel_delay.py
usr/share/openvswitch/scripts/usdt/kernel_delay.rst
usr/share/openvswitch/scripts/usdt/reval_monitor.py
usr/share/openvswitch/scripts/usdt/upcall_cost.py
usr/share/openvswitch/scripts/usdt/upcall_monitor.py
usr/share/openvswitch/vswitch.ovsschema
usr/share/openvswitch/vtep.ovsschema
var/
var/lib/
var/lib/openvswitch/
var/lib/openvswitch/pki/
