# extracted/darkstat-3.0.721-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/darkstat.conf
usr/
usr/bin/
usr/bin/darkstat
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/darkstat.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/darkstat.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/darkstat.conf
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/darkstat.8.gz
