# extracted/dictd-1.13.3-1/files
%FILES%
etc/
etc/default/
etc/default/dictd
etc/dict/
etc/dict/colorit.conf
etc/dict/dict.conf
etc/dict/dictd.conf
etc/dict/site.info
usr/
usr/bin/
usr/bin/colorit
usr/bin/dict
usr/bin/dict_lookup
usr/bin/dictd
usr/bin/dictdplugin-config
usr/bin/dictfmt
usr/bin/dictfmt_index2suffix
usr/bin/dictfmt_index2word
usr/bin/dictl
usr/bin/dictunformat
usr/bin/dictzip
usr/include/
usr/include/dictdplugin.h
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dictd.service
usr/share/
usr/share/doc/
usr/share/doc/dictd/
usr/share/doc/dictd/examples/
usr/share/doc/dictd/examples/dict1.conf
usr/share/doc/dictd/examples/dictd1.conf
usr/share/doc/dictd/examples/dictd2.conf
usr/share/doc/dictd/examples/dictd3.conf
usr/share/doc/dictd/examples/dictd4.conf
usr/share/doc/dictd/examples/dictd_complex.conf
usr/share/doc/dictd/examples/dictd_mime.conf
usr/share/doc/dictd/examples/dictd_plugin_dbi.conf
usr/share/doc/dictd/examples/dictd_popen.conf
usr/share/doc/dictd/examples/dictd_popen.conf.in
usr/share/doc/dictd/examples/dictd_site.txt
usr/share/doc/dictd/examples/dictd_virtual.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/colorit.1.gz
usr/share/man/man1/dict.1.gz
usr/share/man/man1/dict_lookup.1.gz
usr/share/man/man1/dictfmt.1.gz
usr/share/man/man1/dictfmt_index2suffix.1.gz
usr/share/man/man1/dictfmt_index2word.1.gz
usr/share/man/man1/dictl.1.gz
usr/share/man/man1/dictunformat.1.gz
usr/share/man/man1/dictzip.1.gz
usr/share/man/man8/
usr/share/man/man8/dictd.8.gz
