# extracted/pacredir-0.4.7-3/files
%FILES%
etc/
etc/avahi/
etc/avahi/services/
etc/avahi/services/pacserve.service
etc/pacman.d/
etc/pacman.d/pacredir
etc/pacredir.conf
usr/
usr/bin/
usr/bin/pacredir
usr/bin/pacserve
usr/lib/
usr/lib/NetworkManager/
usr/lib/NetworkManager/dispatcher.d/
usr/lib/NetworkManager/dispatcher.d/80-pacredir
usr/lib/dhcpcd/
usr/lib/dhcpcd/dhcpcd-hooks/
usr/lib/dhcpcd/dhcpcd-hooks/80-pacredir
usr/lib/initcpio/
usr/lib/initcpio/hooks/
usr/lib/initcpio/hooks/pacredir
usr/lib/initcpio/install/
usr/lib/initcpio/install/pacredir
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pacredir.service
usr/lib/systemd/system/pacserve.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/pacredir.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/pacredir.conf
usr/share/
usr/share/doc/
usr/share/doc/pacredir/
usr/share/doc/pacredir/COPYING.html
usr/share/doc/pacredir/COPYING.md
usr/share/doc/pacredir/FLOW.html
usr/share/doc/pacredir/FLOW.md
usr/share/doc/pacredir/FLOW/
usr/share/doc/pacredir/FLOW/fail.svg
usr/share/doc/pacredir/FLOW/ok.svg
usr/share/doc/pacredir/README.html
usr/share/doc/pacredir/README.md
