# extracted/ukui-control-center-3.0.4-12/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/com.control.center.qt.systemdbus.conf
etc/dbus-1/system.d/org.ukui.groupmanager.conf
etc/pam.d/
etc/pam.d/control-center
usr/
usr/bin/
usr/bin/changeotheruserpwd
usr/bin/changeuserpwd
usr/bin/checkUserPwd
usr/bin/childCheckpwdwithPAM
usr/bin/group-manager-server
usr/bin/launchSysDbus
usr/bin/ukui-control-center
usr/bin/ukui-control-center-session
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ukui-group-manager.service
usr/lib/ukui-control-center/
usr/lib/ukui-control-center/libabout.so
usr/lib/ukui-control-center/libarea.so
usr/lib/ukui-control-center/libaudio.so
usr/lib/ukui-control-center/libautoboot.so
usr/lib/ukui-control-center/libbackup.so
usr/lib/ukui-control-center/libbluetooth.so
usr/lib/ukui-control-center/libdatetime.so
usr/lib/ukui-control-center/libdefaultapp.so
usr/lib/ukui-control-center/libdesktop.so
usr/lib/ukui-control-center/libdisplay.so
usr/lib/ukui-control-center/libexperienceplan.so
usr/lib/ukui-control-center/libfonts.so
usr/lib/ukui-control-center/libkeyboard.so
usr/lib/ukui-control-center/libmouse.so
usr/lib/ukui-control-center/libnetconnect.so
usr/lib/ukui-control-center/libnetworkaccount.so
usr/lib/ukui-control-center/libnotice.so
usr/lib/ukui-control-center/libpower.so
usr/lib/ukui-control-center/libprinter.so
usr/lib/ukui-control-center/libprojection.so
usr/lib/ukui-control-center/libproxy.so
usr/lib/ukui-control-center/libscreenlock.so
usr/lib/ukui-control-center/libscreensaver.so
usr/lib/ukui-control-center/libsearch.so
usr/lib/ukui-control-center/libshortcut.so
usr/lib/ukui-control-center/libtheme.so
usr/lib/ukui-control-center/libtouchpad.so
usr/lib/ukui-control-center/libtouchscreen.so
usr/lib/ukui-control-center/libupdate.so
usr/lib/ukui-control-center/libuserinfo.so
usr/lib/ukui-control-center/libvino.so
usr/lib/ukui-control-center/libvpn.so
usr/lib/ukui-control-center/libwallpaper.so
usr/share/
usr/share/applications/
usr/share/applications/ukui-control-center.desktop
usr/share/dbus-1/
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.ukui.ukcc.session.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/com.control.center.qt.systemdbus.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/com.control.center.qt.systemdbus.conf
usr/share/glib-2.0/
usr/share/glib-2.0/schemas/
usr/share/glib-2.0/schemas/org.ukui.control-center.desktop.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.experienceplan.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.keybinding.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.keyboard.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.notice.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.panel.plugins.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.personalise.gschema.xml
usr/share/glib-2.0/schemas/org.ukui.control-center.wifi.switch.gschema.xml
usr/share/kylin-user-guide/
usr/share/kylin-user-guide/data/
usr/share/kylin-user-guide/data/guide/
usr/share/kylin-user-guide/data/guide/ukui-control-center/
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/01.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/02.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/03.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/04.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/05.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/06.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/07.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/08.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/09.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/1.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/10.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/11.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/12.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/13.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/14.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/15.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/16.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/17.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/18.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/19.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/2.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/20.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/21.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/22.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/23.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/24.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/25.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/26.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/27.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/28.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/29.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/3.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/30.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/31.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/32.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/33.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/34.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/35.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/36.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/37.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/38.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/39.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/4.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/40.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/41.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/42.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/43.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/44.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/45.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/46.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/47.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/5.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/6.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/7.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/8.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/9.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/about.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account-add.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account-edit.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account-face.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account-pwd.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account-type.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/account.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/add-autoboot.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/add-shortcut.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/area-format.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/area.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/audio.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/autoboot.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/background.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/backup.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/blutooth.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/choicewidget.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/cloudaccount.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/datetime-change.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/datetime-zone.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/datetime.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/default.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/delegate.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/desktop.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/display.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/font.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/historylog.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/keyboard.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/mouse.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/netconnect.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/notice.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/notiend.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/notiffirst.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/power-custom.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/power.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/printer.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/screenlock.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/screensaver.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/security.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/shortcut.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/theme-cursor.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/theme-effect.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/theme.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/touchpad.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/touchscreen.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/tray.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/ukcc.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/update.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/vino.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/image/vpn.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/en_US/index.md
usr/share/kylin-user-guide/data/guide/ukui-control-center/ukui-control-center.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/about.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/accessInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/accessPinInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account-add.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account-edit.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account-face.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account-pwd.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account-type.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/account.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/add-autoboot.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/add-shortcut.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/area-format.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/area.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/audio.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/autoboot.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/background.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/backup.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/blutooth.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/choicewidget.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/cloudaccount.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/datetime-change.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/datetime-zone.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/datetime.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/default.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/delegate.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/desktop.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/display.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/exitInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/font.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/historylog.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/horProjectionInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/keyboard.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/mainInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/mouse.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/nameChangeInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/netconnect.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/notice.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/notiend.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/notiffirst.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/openInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/phoneConnectionInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/phoneOpengInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/phonePinInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/phoneSearchInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/power-custom.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/power.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/printer.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/receiverConnectionInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/screenlock.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/screensaver.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/search.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/security.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/shortcut.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/theme-cursor.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/theme-effect.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/theme.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/touchpad.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/touchscreen.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/tray.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/ukcc.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/unUseInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/unUseInterface2.jpg
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/unUseInterface3.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/update.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/verProjectionInterface.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/vino.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/image/vpn.png
usr/share/kylin-user-guide/data/guide/ukui-control-center/zh_CN/index.md
usr/share/locale/
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/installer-timezones.mo
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/com.control.center.qt.systemdbus.policy
usr/share/polkit-1/actions/org.ukui.groupmanager.policy
usr/share/ukui-control-center/
usr/share/ukui-control-center/shell/
usr/share/ukui-control-center/shell/res/
usr/share/ukui-control-center/shell/res/i18n/
usr/share/ukui-control-center/shell/res/i18n/bo.qm
usr/share/ukui-control-center/shell/res/i18n/bo_CN.qm
usr/share/ukui-control-center/shell/res/i18n/en_US.qm
usr/share/ukui-control-center/shell/res/i18n/fa.qm
usr/share/ukui-control-center/shell/res/i18n/tr.qm
usr/share/ukui-control-center/shell/res/i18n/zh_CN.qm
usr/share/ukui-control-center/shell/res/search.xml
usr/share/ukui/
usr/share/ukui/faces/
usr/share/ukui/faces/1.png
usr/share/ukui/faces/10.png
usr/share/ukui/faces/2.png
usr/share/ukui/faces/3.png
usr/share/ukui/faces/4.png
usr/share/ukui/faces/5.png
usr/share/ukui/faces/6.png
usr/share/ukui/faces/7.png
usr/share/ukui/faces/8.png
usr/share/ukui/faces/9.png
usr/share/ukui/faces/default.png
usr/share/ukui/faces/default.svg
