# extracted/hblock-3.5.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/hblock
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hblock.service
usr/lib/systemd/system/hblock.timer
usr/share/
usr/share/doc/
usr/share/doc/hblock/
usr/share/doc/hblock/README.md
usr/share/doc/hblock/alt-formats/
usr/share/doc/hblock/alt-formats/adblock.txt.sh
usr/share/doc/hblock/alt-formats/dnsmasq.conf.sh
usr/share/doc/hblock/alt-formats/domains.txt.sh
usr/share/doc/hblock/alt-formats/rpz.txt.sh
usr/share/doc/hblock/alt-formats/rpz_wildcards.txt.sh
usr/share/doc/hblock/alt-formats/unbound.conf.sh
usr/share/doc/hblock/alt-formats/windows.ps1.sh
usr/share/licenses/
usr/share/licenses/hblock/
usr/share/licenses/hblock/LICENSE.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hblock.1.gz
