# extracted/switcheroo-control-2.6-4/files
%FILES%
usr/
usr/bin/
usr/bin/switcherooctl
usr/lib/
usr/lib/switcheroo-control
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/switcheroo-control.service
usr/lib/udev/
usr/lib/udev/hwdb.d/
usr/lib/udev/hwdb.d/30-pci-intel-gpu.hwdb
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/net.hadess.SwitcherooControl.conf
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/switcheroo-control/
usr/share/gtk-doc/html/switcheroo-control/gdbus-net.hadess.SwitcherooControl.html
usr/share/gtk-doc/html/switcheroo-control/home.png
usr/share/gtk-doc/html/switcheroo-control/index.html
usr/share/gtk-doc/html/switcheroo-control/ix01.html
usr/share/gtk-doc/html/switcheroo-control/left-insensitive.png
usr/share/gtk-doc/html/switcheroo-control/left.png
usr/share/gtk-doc/html/switcheroo-control/license.html
usr/share/gtk-doc/html/switcheroo-control/ref-dbus.html
usr/share/gtk-doc/html/switcheroo-control/right-insensitive.png
usr/share/gtk-doc/html/switcheroo-control/right.png
usr/share/gtk-doc/html/switcheroo-control/style.css
usr/share/gtk-doc/html/switcheroo-control/switcheroo-control.devhelp2
usr/share/gtk-doc/html/switcheroo-control/switcherooctl.html
usr/share/gtk-doc/html/switcheroo-control/tools.html
usr/share/gtk-doc/html/switcheroo-control/up-insensitive.png
usr/share/gtk-doc/html/switcheroo-control/up.png
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/switcherooctl.1.gz
