# extracted/dnscrypt-wrapper-0.4.2-4/files
%FILES%
etc/
etc/default/
etc/default/dnscrypt-wrapper
usr/
usr/bin/
usr/bin/dnscrypt-wrapper
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dnscrypt-wrapper.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/dnscrypt-wrapper.conf
usr/share/
usr/share/licenses/
usr/share/licenses/dnscrypt-wrapper/
usr/share/licenses/dnscrypt-wrapper/COPYING
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/dnscrypt-wrapper.8.gz
