# extracted/gotosocial-0.17.3-1/files
%FILES%
etc/
etc/gotosocial/
etc/gotosocial/config.yaml
etc/gotosocial/template/
etc/gotosocial/template/404.tmpl
etc/gotosocial/template/about.tmpl
etc/gotosocial/template/authorize.tmpl
etc/gotosocial/template/confirm_email.tmpl
etc/gotosocial/template/confirmed_email.tmpl
etc/gotosocial/template/domain-blocklist.tmpl
etc/gotosocial/template/email_confirm.tmpl
etc/gotosocial/template/email_new_report.tmpl
etc/gotosocial/template/email_new_signup.tmpl
etc/gotosocial/template/email_report_closed.tmpl
etc/gotosocial/template/email_reset.tmpl
etc/gotosocial/template/email_signup_approved.tmpl
etc/gotosocial/template/email_signup_rejected.tmpl
etc/gotosocial/template/email_test.tmpl
etc/gotosocial/template/error.tmpl
etc/gotosocial/template/finalize.tmpl
etc/gotosocial/template/index.tmpl
etc/gotosocial/template/index_apps.tmpl
etc/gotosocial/template/index_register.tmpl
etc/gotosocial/template/index_what_is_this.tmpl
etc/gotosocial/template/oob.tmpl
etc/gotosocial/template/page.tmpl
etc/gotosocial/template/page_footer.tmpl
etc/gotosocial/template/page_header.tmpl
etc/gotosocial/template/page_ogmeta.tmpl
etc/gotosocial/template/page_stylesheets.tmpl
etc/gotosocial/template/profile.tmpl
etc/gotosocial/template/profile_fields.tmpl
etc/gotosocial/template/settings.tmpl
etc/gotosocial/template/sign-in.tmpl
etc/gotosocial/template/sign-up.tmpl
etc/gotosocial/template/signed-up.tmpl
etc/gotosocial/template/status.tmpl
etc/gotosocial/template/status_attachments.tmpl
etc/gotosocial/template/status_attributes.tmpl
etc/gotosocial/template/status_header.tmpl
etc/gotosocial/template/status_info.tmpl
etc/gotosocial/template/status_poll.tmpl
etc/gotosocial/template/tag.tmpl
etc/gotosocial/template/thread.tmpl
usr/
usr/bin/
usr/bin/gotosocial
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/gotosocial.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/gotosocial.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/gotosocial.conf
usr/share/
usr/share/gotosocial/
usr/share/gotosocial/Fork-Awesome/
usr/share/gotosocial/Fork-Awesome/CONTRIBUTORS.md
usr/share/gotosocial/Fork-Awesome/LICENSES
usr/share/gotosocial/Fork-Awesome/css/
usr/share/gotosocial/Fork-Awesome/css/fork-awesome.css
usr/share/gotosocial/Fork-Awesome/css/fork-awesome.min.css
usr/share/gotosocial/Fork-Awesome/css/fork-awesome.min.css.map
usr/share/gotosocial/Fork-Awesome/css/v5-compat.css
usr/share/gotosocial/Fork-Awesome/css/v5-compat.min.css
usr/share/gotosocial/Fork-Awesome/css/v5-compat.min.css.map
usr/share/gotosocial/Fork-Awesome/fonts/
usr/share/gotosocial/Fork-Awesome/fonts/forkawesome-webfont.eot
usr/share/gotosocial/Fork-Awesome/fonts/forkawesome-webfont.svg
usr/share/gotosocial/Fork-Awesome/fonts/forkawesome-webfont.ttf
usr/share/gotosocial/Fork-Awesome/fonts/forkawesome-webfont.woff
usr/share/gotosocial/Fork-Awesome/fonts/forkawesome-webfont.woff2
usr/share/gotosocial/LICENSE
usr/share/gotosocial/activitypub_light.svg
usr/share/gotosocial/default_avatars/
usr/share/gotosocial/default_avatars/GoToSocial_icon1.webp
usr/share/gotosocial/default_avatars/GoToSocial_icon2.webp
usr/share/gotosocial/default_avatars/GoToSocial_icon3.webp
usr/share/gotosocial/default_avatars/GoToSocial_icon4.webp
usr/share/gotosocial/default_avatars/GoToSocial_icon5.webp
usr/share/gotosocial/default_avatars/GoToSocial_icon6.webp
usr/share/gotosocial/default_header.webp
usr/share/gotosocial/dist/
usr/share/gotosocial/dist/_colors.css
usr/share/gotosocial/dist/about.css
usr/share/gotosocial/dist/base.css
usr/share/gotosocial/dist/frontend.js
usr/share/gotosocial/dist/index.css
usr/share/gotosocial/dist/page.css
usr/share/gotosocial/dist/prism.css
usr/share/gotosocial/dist/profile.css
usr/share/gotosocial/dist/settings-style.css
usr/share/gotosocial/dist/settings.js
usr/share/gotosocial/dist/status.css
usr/share/gotosocial/dist/tag.css
usr/share/gotosocial/dist/thread.css
usr/share/gotosocial/feditext.svg
usr/share/gotosocial/fonts/
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Bold-102.woff
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Bold-102a.woff2
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-BoldItalic-102.woff
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-BoldItalic-102a.woff2
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Italic-102.woff
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Italic-102a.woff2
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Regular-102.woff
usr/share/gotosocial/fonts/Atkinson-Hyperlegible-Regular-102a.woff2
usr/share/gotosocial/fonts/OFL.txt
usr/share/gotosocial/fonts/noto-sans-mono-v30-latin-regular.woff2
usr/share/gotosocial/fonts/noto-sans-v27-latin-700.woff
usr/share/gotosocial/fonts/noto-sans-v27-latin-700.woff2
usr/share/gotosocial/fonts/noto-sans-v27-latin-regular.woff
usr/share/gotosocial/fonts/noto-sans-v27-latin-regular.woff2
usr/share/gotosocial/logo.png
usr/share/gotosocial/logo.svg
usr/share/gotosocial/logo.webp
usr/share/gotosocial/mastodon.svg
usr/share/gotosocial/plyr.svg
usr/share/gotosocial/swagger.yaml
usr/share/gotosocial/themes/
usr/share/gotosocial/themes/blurple-dark.css
usr/share/gotosocial/themes/blurple-light.css
usr/share/gotosocial/themes/brutalist-dark.css
usr/share/gotosocial/themes/brutalist.css
usr/share/gotosocial/themes/ecks-pee.css
usr/share/gotosocial/themes/midnight-trip.css
usr/share/gotosocial/themes/moonlight-hunt.css
usr/share/gotosocial/themes/rain-forest.css
usr/share/gotosocial/themes/soft.css
usr/share/gotosocial/themes/solarized-dark.css
usr/share/gotosocial/themes/solarized-light.css
usr/share/gotosocial/themes/sunset-light.css
usr/share/gotosocial/tusky.svg
