# extracted/slim-1:1.4.0-1/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/slim
etc/pam.d/
etc/pam.d/slim
etc/pam.d/slimlock
etc/slim.conf
etc/slimlock.conf
usr/
usr/bin/
usr/bin/slim
usr/bin/slimlock
usr/lib/
usr/lib/libslim.so
usr/lib/libslim.so.1.4.0
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/slim.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/slim.1.gz
usr/share/man/man1/slimlock.1.gz
usr/share/slim/
usr/share/slim/themes/
usr/share/slim/themes/default/
usr/share/slim/themes/default/background.jpg
usr/share/slim/themes/default/panel.png
usr/share/slim/themes/default/slim.theme
usr/share/slim/themes/original/
usr/share/slim/themes/original/background.jpg
usr/share/slim/themes/original/panel.png
usr/share/slim/themes/original/slim.theme
