# extracted/coturn-4.6.2-2/files
%FILES%
etc/
etc/turnserver/
etc/turnserver/turnserver.conf
usr/
usr/bin/
usr/bin/turnadmin
usr/bin/turnserver
usr/bin/turnutils_natdiscovery
usr/bin/turnutils_oauth
usr/bin/turnutils_peer
usr/bin/turnutils_stunclient
usr/bin/turnutils_uclient
usr/include/
usr/include/turn/
usr/include/turn/client/
usr/include/turn/client/TurnMsgLib.h
usr/include/turn/client/ns_turn_ioaddr.h
usr/include/turn/client/ns_turn_msg.h
usr/include/turn/client/ns_turn_msg_addr.h
usr/include/turn/client/ns_turn_msg_defs.h
usr/include/turn/client/ns_turn_msg_defs_experimental.h
usr/include/turn/ns_turn_defs.h
usr/lib/
usr/lib/libturnclient.a
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/coturn.service
usr/lib/systemd/system/turnserver.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/turnserver.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/turnserver.conf
usr/share/
usr/share/doc/
usr/share/doc/turnserver/
usr/share/doc/turnserver/INSTALL
usr/share/doc/turnserver/LICENSE
usr/share/doc/turnserver/README.turnadmin
usr/share/doc/turnserver/README.turnserver
usr/share/doc/turnserver/README.turnutils
usr/share/doc/turnserver/postinstall.txt
usr/share/doc/turnserver/schema.mongo.sh
usr/share/doc/turnserver/schema.sql
usr/share/doc/turnserver/schema.stats.redis
usr/share/doc/turnserver/schema.userdb.redis
usr/share/licenses/
usr/share/licenses/coturn/
usr/share/licenses/coturn/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/coturn.1.gz
usr/share/man/man1/turnadmin.1.gz
usr/share/man/man1/turnserver.1.gz
usr/share/man/man1/turnutils.1.gz
usr/share/man/man1/turnutils_natdiscovery.1.gz
usr/share/man/man1/turnutils_oauth.1.gz
usr/share/man/man1/turnutils_peer.1.gz
usr/share/man/man1/turnutils_stunclient.1.gz
usr/share/man/man1/turnutils_uclient.1.gz
usr/share/turnserver/
usr/share/turnserver/examples/
usr/share/turnserver/examples/etc/
usr/share/turnserver/examples/etc/cacert.pem
usr/share/turnserver/examples/etc/coturn.service
usr/share/turnserver/examples/etc/turn_client_cert.pem
usr/share/turnserver/examples/etc/turn_client_pkey.pem
usr/share/turnserver/examples/etc/turn_server_cert.pem
usr/share/turnserver/examples/etc/turn_server_pkey.pem
usr/share/turnserver/examples/etc/turnserver.conf
usr/share/turnserver/examples/scripts/
usr/share/turnserver/examples/scripts/basic/
usr/share/turnserver/examples/scripts/basic/dos_attack.sh
usr/share/turnserver/examples/scripts/basic/relay.sh
usr/share/turnserver/examples/scripts/basic/tcp_client.sh
usr/share/turnserver/examples/scripts/basic/tcp_client_c2c_tcp_relay.sh
usr/share/turnserver/examples/scripts/basic/udp_c2c_client.sh
usr/share/turnserver/examples/scripts/basic/udp_client.sh
usr/share/turnserver/examples/scripts/loadbalance/
usr/share/turnserver/examples/scripts/loadbalance/master_relay.sh
usr/share/turnserver/examples/scripts/loadbalance/slave_relay_1.sh
usr/share/turnserver/examples/scripts/loadbalance/slave_relay_2.sh
usr/share/turnserver/examples/scripts/loadbalance/tcp_c2c_tcp_relay.sh
usr/share/turnserver/examples/scripts/loadbalance/udp_c2c.sh
usr/share/turnserver/examples/scripts/longtermsecure/
usr/share/turnserver/examples/scripts/longtermsecure/secure_dos_attack.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_dtls_client.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_dtls_client_cert.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_relay.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_relay_cert.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_sctp_client.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_tcp_client.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_tcp_client_c2c_tcp_relay.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_tls_client.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_tls_client_c2c_tcp_relay.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_tls_client_cert.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_udp_c2c.sh
usr/share/turnserver/examples/scripts/longtermsecure/secure_udp_client.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_mongo.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_mysql.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_mysql_ssl.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_psql.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_redis.sh
usr/share/turnserver/examples/scripts/longtermsecuredb/secure_relay_with_db_sqlite.sh
usr/share/turnserver/examples/scripts/mobile/
usr/share/turnserver/examples/scripts/mobile/mobile_dtls_client.sh
usr/share/turnserver/examples/scripts/mobile/mobile_relay.sh
usr/share/turnserver/examples/scripts/mobile/mobile_tcp_client.sh
usr/share/turnserver/examples/scripts/mobile/mobile_tls_client_c2c_tcp_relay.sh
usr/share/turnserver/examples/scripts/mobile/mobile_udp_client.sh
usr/share/turnserver/examples/scripts/oauth.sh
usr/share/turnserver/examples/scripts/pack.sh
usr/share/turnserver/examples/scripts/peer.sh
usr/share/turnserver/examples/scripts/readme.txt
usr/share/turnserver/examples/scripts/restapi/
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret.sh
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret_with_db_mongo.sh
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret_with_db_mysql.sh
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret_with_db_psql.sh
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret_with_db_redis.sh
usr/share/turnserver/examples/scripts/restapi/secure_relay_secret_with_db_sqlite.sh
usr/share/turnserver/examples/scripts/restapi/secure_udp_client_with_secret.sh
usr/share/turnserver/examples/scripts/restapi/shared_secret_maintainer.pl
usr/share/turnserver/examples/scripts/selfloadbalance/
usr/share/turnserver/examples/scripts/selfloadbalance/secure_dos_attack.sh
usr/share/turnserver/examples/scripts/selfloadbalance/secure_relay.sh
usr/share/turnserver/schema.mongo.sh
usr/share/turnserver/schema.sql
usr/share/turnserver/schema.stats.redis
usr/share/turnserver/schema.userdb.redis
usr/share/turnserver/testmongosetup.sh
usr/share/turnserver/testredisdbsetup.sh
usr/share/turnserver/testsqldbsetup.sql
var/
var/db/
var/db/turndb
