# extracted/ovenmediaengine-0.17.2-1/files
%FILES%
etc/
etc/ovenmediaengine/
etc/ovenmediaengine/Logger.xml
etc/ovenmediaengine/Server.xml
usr/
usr/bin/
usr/bin/OvenMediaEngine
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ovenmediaengine.service
usr/share/
usr/share/doc/
usr/share/doc/ovenmediaengine/
usr/share/doc/ovenmediaengine/README.md
usr/share/doc/ovenmediaengine/SUMMARY.md
usr/share/doc/ovenmediaengine/access-control/
usr/share/doc/ovenmediaengine/access-control/README.md
usr/share/doc/ovenmediaengine/access-control/admission-webhooks.md
usr/share/doc/ovenmediaengine/access-control/signedpolicy.md
usr/share/doc/ovenmediaengine/alert.md
usr/share/doc/ovenmediaengine/conf_examples/
usr/share/doc/ovenmediaengine/conf_examples/DrmInfo.xml
usr/share/doc/ovenmediaengine/conf_examples/Edge.xml
usr/share/doc/ovenmediaengine/conf_examples/Logger.xml
usr/share/doc/ovenmediaengine/conf_examples/Origin.xml
usr/share/doc/ovenmediaengine/conf_examples/Server.xml
usr/share/doc/ovenmediaengine/configuration/
usr/share/doc/ovenmediaengine/configuration/README.md
usr/share/doc/ovenmediaengine/configuration/ipv6.md
usr/share/doc/ovenmediaengine/configuration/tls-encryption.md
usr/share/doc/ovenmediaengine/crossdomains.md
usr/share/doc/ovenmediaengine/getting-started/
usr/share/doc/ovenmediaengine/getting-started/README.md
usr/share/doc/ovenmediaengine/getting-started/getting-started-with-docker.md
usr/share/doc/ovenmediaengine/getting-started/getting-started-with-ome-docker-launcher.md
usr/share/doc/ovenmediaengine/live-source/
usr/share/doc/ovenmediaengine/live-source/README.md
usr/share/doc/ovenmediaengine/live-source/mpeg-2-ts-beta.md
usr/share/doc/ovenmediaengine/live-source/multiplex-channel.md
usr/share/doc/ovenmediaengine/live-source/rtmp.md
usr/share/doc/ovenmediaengine/live-source/rtsp-pull.md
usr/share/doc/ovenmediaengine/live-source/scheduled-channel.md
usr/share/doc/ovenmediaengine/live-source/srt.md
usr/share/doc/ovenmediaengine/live-source/webrtc.md
usr/share/doc/ovenmediaengine/logs-and-statistics.md
usr/share/doc/ovenmediaengine/origin-edge-clustering.md
usr/share/doc/ovenmediaengine/p2p-delivery.md
usr/share/doc/ovenmediaengine/performance-tuning.md
usr/share/doc/ovenmediaengine/push-publishing.md
usr/share/doc/ovenmediaengine/quick-start/
usr/share/doc/ovenmediaengine/quick-start/README.md
usr/share/doc/ovenmediaengine/quick-start/test-player.md
usr/share/doc/ovenmediaengine/recording.md
usr/share/doc/ovenmediaengine/rest-api/
usr/share/doc/ovenmediaengine/rest-api/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/
usr/share/doc/ovenmediaengine/rest-api/v1/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/statistics/
usr/share/doc/ovenmediaengine/rest-api/v1/statistics/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/statistics/current.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/output-profile.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/push.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/recording.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/scheduledchannel-api-1.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/scheduledchannel-api.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/stream/
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/stream/README.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/stream/conclude-hls-live.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/stream/hls-dump.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/application/stream/send-event.md
usr/share/doc/ovenmediaengine/rest-api/v1/virtualhost/reload-certificate.md
usr/share/doc/ovenmediaengine/streaming/
usr/share/doc/ovenmediaengine/streaming/README.md
usr/share/doc/ovenmediaengine/streaming/hls.md
usr/share/doc/ovenmediaengine/streaming/low-latency-hls.md
usr/share/doc/ovenmediaengine/streaming/webrtc-publishing.md
usr/share/doc/ovenmediaengine/thumbnail.md
usr/share/doc/ovenmediaengine/transcoding/
usr/share/doc/ovenmediaengine/transcoding/README.md
usr/share/doc/ovenmediaengine/transcoding/gpu-usage.md
usr/share/doc/ovenmediaengine/transcoding/transcodewebhook.md
usr/share/doc/ovenmediaengine/troubleshooting.md
usr/share/licenses/
usr/share/licenses/ovenmediaengine/
usr/share/licenses/ovenmediaengine/LICENSE
