# extracted/fetchmail-6.5.1-1/files
%FILES%
usr/
usr/bin/
usr/bin/fetchmail
usr/bin/fetchmailconf
usr/lib/
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/__pycache__/
usr/lib/python3.13/site-packages/__pycache__/fetchmailconf.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/__pycache__/fetchmailconf.cpython-313.pyc
usr/lib/python3.13/site-packages/fetchmailconf.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/fetchmail.service
usr/lib/systemd/system/fetchmail@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/fetchmail.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/fetchmail.conf
usr/share/
usr/share/locale/
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/fetchmail.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/fetchmail.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/fetchmail.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/fetchmail.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/fetchmail.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/fetchmail.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/fetchmail.mo
usr/share/locale/ro/
usr/share/locale/ro/LC_MESSAGES/
usr/share/locale/ro/LC_MESSAGES/fetchmail.mo
usr/share/locale/sq/
usr/share/locale/sq/LC_MESSAGES/
usr/share/locale/sq/LC_MESSAGES/fetchmail.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/fetchmail.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/fetchmail.1.gz
usr/share/man/man1/fetchmailconf.1.gz
