# extracted/xf86-input-wacom-1.2.2-1/files
%FILES%
usr/
usr/bin/
usr/bin/isdv4-serial-debugger
usr/bin/isdv4-serial-inputattach
usr/bin/xsetwacom
usr/include/
usr/include/xorg/
usr/include/xorg/Xwacom.h
usr/include/xorg/isdv4.h
usr/include/xorg/wacom-properties.h
usr/include/xorg/wacom-util.h
usr/lib/
usr/lib/pkgconfig/
usr/lib/pkgconfig/xorg-wacom.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wacom-inputattach@.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/wacom.rules
usr/lib/xorg/
usr/lib/xorg/modules/
usr/lib/xorg/modules/input/
usr/lib/xorg/modules/input/wacom_drv.so
usr/share/
usr/share/X11/
usr/share/X11/xorg.conf.d/
usr/share/X11/xorg.conf.d/70-wacom.conf
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/xsetwacom.1.gz
usr/share/man/man4/
usr/share/man/man4/wacom.4.gz
