# extracted/sndio-1.10.0-1/files
%FILES%
etc/
etc/default/
etc/default/sndiod
usr/
usr/bin/
usr/bin/aucat
usr/bin/midicat
usr/bin/sndioctl
usr/bin/sndiod
usr/include/
usr/include/sndio.h
usr/lib/
usr/lib/libsndio.so
usr/lib/libsndio.so.7
usr/lib/libsndio.so.7.0
usr/lib/libsndio.so.7.3
usr/lib/pkgconfig/
usr/lib/pkgconfig/sndio.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sndiod.service
usr/share/
usr/share/licenses/
usr/share/licenses/sndio/
usr/share/licenses/sndio/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/aucat.1.gz
usr/share/man/man1/midicat.1.gz
usr/share/man/man1/sndioctl.1.gz
usr/share/man/man3/
usr/share/man/man3/mio_close.3.gz
usr/share/man/man3/mio_eof.3.gz
usr/share/man/man3/mio_nfds.3.gz
usr/share/man/man3/mio_open.3.gz
usr/share/man/man3/mio_pollfd.3.gz
usr/share/man/man3/mio_read.3.gz
usr/share/man/man3/mio_revents.3.gz
usr/share/man/man3/mio_write.3.gz
usr/share/man/man3/sio_close.3.gz
usr/share/man/man3/sio_eof.3.gz
usr/share/man/man3/sio_getcap.3.gz
usr/share/man/man3/sio_getpar.3.gz
usr/share/man/man3/sio_initpar.3.gz
usr/share/man/man3/sio_nfds.3.gz
usr/share/man/man3/sio_onmove.3.gz
usr/share/man/man3/sio_onvol.3.gz
usr/share/man/man3/sio_open.3.gz
usr/share/man/man3/sio_pollfd.3.gz
usr/share/man/man3/sio_read.3.gz
usr/share/man/man3/sio_revents.3.gz
usr/share/man/man3/sio_setpar.3.gz
usr/share/man/man3/sio_setvol.3.gz
usr/share/man/man3/sio_start.3.gz
usr/share/man/man3/sio_stop.3.gz
usr/share/man/man3/sio_write.3.gz
usr/share/man/man3/sioctl_close.3.gz
usr/share/man/man3/sioctl_eof.3.gz
usr/share/man/man3/sioctl_nfds.3.gz
usr/share/man/man3/sioctl_ondesc.3.gz
usr/share/man/man3/sioctl_onval.3.gz
usr/share/man/man3/sioctl_open.3.gz
usr/share/man/man3/sioctl_pollfd.3.gz
usr/share/man/man3/sioctl_revents.3.gz
usr/share/man/man3/sioctl_setval.3.gz
usr/share/man/man7/
usr/share/man/man7/sndio.7.gz
usr/share/man/man8/
usr/share/man/man8/sndiod.8.gz
