# extracted/dnscrypt-proxy-2.1.5-4/files
%FILES%
etc/
etc/dnscrypt-proxy/
etc/dnscrypt-proxy/allowed-ips.txt
etc/dnscrypt-proxy/allowed-names.txt
etc/dnscrypt-proxy/blocked-ips.txt
etc/dnscrypt-proxy/blocked-names.txt
etc/dnscrypt-proxy/captive-portals.txt
etc/dnscrypt-proxy/cloaking-rules.txt
etc/dnscrypt-proxy/dnscrypt-proxy.toml
etc/dnscrypt-proxy/forwarding-rules.txt
usr/
usr/bin/
usr/bin/dnscrypt-proxy
usr/bin/generate-domains-blocklist
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dnscrypt-proxy.service
usr/lib/systemd/system/dnscrypt-proxy.socket
usr/share/
usr/share/dnscrypt-proxy/
usr/share/dnscrypt-proxy/utils/
usr/share/dnscrypt-proxy/utils/generate-domains-blocklist/
usr/share/dnscrypt-proxy/utils/generate-domains-blocklist/domains-allowlist.txt
usr/share/dnscrypt-proxy/utils/generate-domains-blocklist/domains-blocklist-local-additions.txt
usr/share/dnscrypt-proxy/utils/generate-domains-blocklist/domains-blocklist.conf
usr/share/dnscrypt-proxy/utils/generate-domains-blocklist/domains-time-restricted.txt
usr/share/doc/
usr/share/doc/dnscrypt-proxy/
usr/share/doc/dnscrypt-proxy/ChangeLog
usr/share/doc/dnscrypt-proxy/README.md
usr/share/licenses/
usr/share/licenses/dnscrypt-proxy/
usr/share/licenses/dnscrypt-proxy/LICENSE
