# extracted/spice-vdagent-0.22.1-4/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/spice-vdagentd
etc/xdg/
etc/xdg/autostart/
etc/xdg/autostart/spice-vdagent.desktop
usr/
usr/bin/
usr/bin/spice-vdagent
usr/bin/spice-vdagentd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/spice-vdagentd.service
usr/lib/systemd/system/spice-vdagentd.socket
usr/lib/systemd/user/
usr/lib/systemd/user/spice-vdagent.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/spice-vdagentd.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/70-spice-vdagentd.rules
usr/share/
usr/share/gdm/
usr/share/gdm/autostart/
usr/share/gdm/autostart/LoginWindow/
usr/share/gdm/autostart/LoginWindow/spice-vdagent.desktop
usr/share/gdm/greeter/
usr/share/gdm/greeter/autostart/
usr/share/gdm/greeter/autostart/spice-vdagent.desktop
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/spice-vdagent.1.gz
usr/share/man/man1/spice-vdagentd.1.gz
