# extracted/kafka-3.9.0-1/files
%FILES%
etc/
etc/kafka/
etc/kafka/connect-console-sink.properties
etc/kafka/connect-console-source.properties
etc/kafka/connect-distributed.properties
etc/kafka/connect-file-sink.properties
etc/kafka/connect-file-source.properties
etc/kafka/connect-log4j.properties
etc/kafka/connect-mirror-maker.properties
etc/kafka/connect-standalone.properties
etc/kafka/consumer.properties
etc/kafka/kraft
etc/kafka/log4j.properties
etc/kafka/producer.properties
etc/kafka/server.properties
etc/kafka/tools-log4j.properties
etc/kafka/trogdor.conf
etc/kafka/zookeeper.properties
usr/
usr/bin/
usr/bin/connect-distributed.sh
usr/bin/connect-mirror-maker.sh
usr/bin/connect-plugin-path.sh
usr/bin/connect-standalone.sh
usr/bin/kafka-acls.sh
usr/bin/kafka-broker-api-versions.sh
usr/bin/kafka-client-metrics.sh
usr/bin/kafka-cluster.sh
usr/bin/kafka-configs.sh
usr/bin/kafka-console-consumer.sh
usr/bin/kafka-console-producer.sh
usr/bin/kafka-consumer-groups.sh
usr/bin/kafka-consumer-perf-test.sh
usr/bin/kafka-delegation-tokens.sh
usr/bin/kafka-delete-records.sh
usr/bin/kafka-dump-log.sh
usr/bin/kafka-e2e-latency.sh
usr/bin/kafka-features.sh
usr/bin/kafka-get-offsets.sh
usr/bin/kafka-jmx.sh
usr/bin/kafka-leader-election.sh
usr/bin/kafka-log-dirs.sh
usr/bin/kafka-metadata-quorum.sh
usr/bin/kafka-metadata-shell.sh
usr/bin/kafka-mirror-maker.sh
usr/bin/kafka-producer-perf-test.sh
usr/bin/kafka-reassign-partitions.sh
usr/bin/kafka-replica-verification.sh
usr/bin/kafka-run-class.sh
usr/bin/kafka-server-start.sh
usr/bin/kafka-server-stop.sh
usr/bin/kafka-storage.sh
usr/bin/kafka-streams-application-reset.sh
usr/bin/kafka-topics.sh
usr/bin/kafka-transactions.sh
usr/bin/kafka-verifiable-consumer.sh
usr/bin/kafka-verifiable-producer.sh
usr/bin/trogdor.sh
usr/bin/zookeeper-security-migration.sh
usr/bin/zookeeper-server-start.sh
usr/bin/zookeeper-server-stop.sh
usr/bin/zookeeper-shell.sh
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/kafka.service
usr/lib/systemd/system/zookeeper.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/kafka.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/kafka.conf
usr/share/
usr/share/doc/
usr/share/doc/kafka/
usr/share/doc/kafka/LICENSE
usr/share/doc/kafka/NOTICE
usr/share/kafka/
usr/share/kafka/bin/
usr/share/kafka/bin/connect-distributed.sh
usr/share/kafka/bin/connect-mirror-maker.sh
usr/share/kafka/bin/connect-plugin-path.sh
usr/share/kafka/bin/connect-standalone.sh
usr/share/kafka/bin/kafka-acls.sh
usr/share/kafka/bin/kafka-broker-api-versions.sh
usr/share/kafka/bin/kafka-client-metrics.sh
usr/share/kafka/bin/kafka-cluster.sh
usr/share/kafka/bin/kafka-configs.sh
usr/share/kafka/bin/kafka-console-consumer.sh
usr/share/kafka/bin/kafka-console-producer.sh
usr/share/kafka/bin/kafka-consumer-groups.sh
usr/share/kafka/bin/kafka-consumer-perf-test.sh
usr/share/kafka/bin/kafka-delegation-tokens.sh
usr/share/kafka/bin/kafka-delete-records.sh
usr/share/kafka/bin/kafka-dump-log.sh
usr/share/kafka/bin/kafka-e2e-latency.sh
usr/share/kafka/bin/kafka-features.sh
usr/share/kafka/bin/kafka-get-offsets.sh
usr/share/kafka/bin/kafka-jmx.sh
usr/share/kafka/bin/kafka-leader-election.sh
usr/share/kafka/bin/kafka-log-dirs.sh
usr/share/kafka/bin/kafka-metadata-quorum.sh
usr/share/kafka/bin/kafka-metadata-shell.sh
usr/share/kafka/bin/kafka-mirror-maker.sh
usr/share/kafka/bin/kafka-producer-perf-test.sh
usr/share/kafka/bin/kafka-reassign-partitions.sh
usr/share/kafka/bin/kafka-replica-verification.sh
usr/share/kafka/bin/kafka-run-class.sh
usr/share/kafka/bin/kafka-run-class.sh.orig
usr/share/kafka/bin/kafka-server-start.sh
usr/share/kafka/bin/kafka-server-stop.sh
usr/share/kafka/bin/kafka-storage.sh
usr/share/kafka/bin/kafka-streams-application-reset.sh
usr/share/kafka/bin/kafka-topics.sh
usr/share/kafka/bin/kafka-transactions.sh
usr/share/kafka/bin/kafka-verifiable-consumer.sh
usr/share/kafka/bin/kafka-verifiable-producer.sh
usr/share/kafka/bin/trogdor.sh
usr/share/kafka/bin/zookeeper-security-migration.sh
usr/share/kafka/bin/zookeeper-server-start.sh
usr/share/kafka/bin/zookeeper-server-stop.sh
usr/share/kafka/bin/zookeeper-shell.sh
usr/share/kafka/config/
usr/share/kafka/config/connect-console-sink.properties
usr/share/kafka/config/connect-console-source.properties
usr/share/kafka/config/connect-distributed.properties
usr/share/kafka/config/connect-file-sink.properties
usr/share/kafka/config/connect-file-source.properties
usr/share/kafka/config/connect-log4j.properties
usr/share/kafka/config/connect-mirror-maker.properties
usr/share/kafka/config/connect-standalone.properties
usr/share/kafka/config/consumer.properties
usr/share/kafka/config/kraft/
usr/share/kafka/config/kraft/broker.properties
usr/share/kafka/config/kraft/controller.properties
usr/share/kafka/config/kraft/reconfig-server.properties
usr/share/kafka/config/kraft/server.properties
usr/share/kafka/config/log4j.properties
usr/share/kafka/config/producer.properties
usr/share/kafka/config/server.properties
usr/share/kafka/config/tools-log4j.properties
usr/share/kafka/config/trogdor.conf
usr/share/kafka/config/zookeeper.properties
usr/share/kafka/libs/
usr/share/kafka/libs/activation-1.1.1.jar
usr/share/kafka/libs/aopalliance-repackaged-2.6.1.jar
usr/share/kafka/libs/argparse4j-0.7.0.jar
usr/share/kafka/libs/audience-annotations-0.12.0.jar
usr/share/kafka/libs/caffeine-2.9.3.jar
usr/share/kafka/libs/commons-beanutils-1.9.4.jar
usr/share/kafka/libs/commons-cli-1.4.jar
usr/share/kafka/libs/commons-collections-3.2.2.jar
usr/share/kafka/libs/commons-digester-2.1.jar
usr/share/kafka/libs/commons-io-2.14.0.jar
usr/share/kafka/libs/commons-lang3-3.12.0.jar
usr/share/kafka/libs/commons-logging-1.2.jar
usr/share/kafka/libs/commons-validator-1.7.jar
usr/share/kafka/libs/connect-api-3.9.0.jar
usr/share/kafka/libs/connect-basic-auth-extension-3.9.0.jar
usr/share/kafka/libs/connect-file-3.9.0.jar
usr/share/kafka/libs/connect-json-3.9.0.jar
usr/share/kafka/libs/connect-mirror-3.9.0.jar
usr/share/kafka/libs/connect-mirror-client-3.9.0.jar
usr/share/kafka/libs/connect-runtime-3.9.0.jar
usr/share/kafka/libs/connect-transforms-3.9.0.jar
usr/share/kafka/libs/error_prone_annotations-2.10.0.jar
usr/share/kafka/libs/hk2-api-2.6.1.jar
usr/share/kafka/libs/hk2-locator-2.6.1.jar
usr/share/kafka/libs/hk2-utils-2.6.1.jar
usr/share/kafka/libs/jackson-annotations-2.16.2.jar
usr/share/kafka/libs/jackson-core-2.16.2.jar
usr/share/kafka/libs/jackson-databind-2.16.2.jar
usr/share/kafka/libs/jackson-dataformat-csv-2.16.2.jar
usr/share/kafka/libs/jackson-datatype-jdk8-2.16.2.jar
usr/share/kafka/libs/jackson-jaxrs-base-2.16.2.jar
usr/share/kafka/libs/jackson-jaxrs-json-provider-2.16.2.jar
usr/share/kafka/libs/jackson-module-afterburner-2.16.2.jar
usr/share/kafka/libs/jackson-module-jaxb-annotations-2.16.2.jar
usr/share/kafka/libs/jackson-module-scala_2.13-2.16.2.jar
usr/share/kafka/libs/jakarta.activation-api-1.2.2.jar
usr/share/kafka/libs/jakarta.annotation-api-1.3.5.jar
usr/share/kafka/libs/jakarta.inject-2.6.1.jar
usr/share/kafka/libs/jakarta.validation-api-2.0.2.jar
usr/share/kafka/libs/jakarta.ws.rs-api-2.1.6.jar
usr/share/kafka/libs/jakarta.xml.bind-api-2.3.3.jar
usr/share/kafka/libs/javassist-3.29.2-GA.jar
usr/share/kafka/libs/javax.activation-api-1.2.0.jar
usr/share/kafka/libs/javax.annotation-api-1.3.2.jar
usr/share/kafka/libs/javax.servlet-api-3.1.0.jar
usr/share/kafka/libs/javax.ws.rs-api-2.1.1.jar
usr/share/kafka/libs/jaxb-api-2.3.1.jar
usr/share/kafka/libs/jersey-client-2.39.1.jar
usr/share/kafka/libs/jersey-common-2.39.1.jar
usr/share/kafka/libs/jersey-container-servlet-2.39.1.jar
usr/share/kafka/libs/jersey-container-servlet-core-2.39.1.jar
usr/share/kafka/libs/jersey-hk2-2.39.1.jar
usr/share/kafka/libs/jersey-server-2.39.1.jar
usr/share/kafka/libs/jetty-client-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-continuation-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-http-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-io-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-security-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-server-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-servlet-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-servlets-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-util-9.4.56.v20240826.jar
usr/share/kafka/libs/jetty-util-ajax-9.4.56.v20240826.jar
usr/share/kafka/libs/jline-3.25.1.jar
usr/share/kafka/libs/jopt-simple-5.0.4.jar
usr/share/kafka/libs/jose4j-0.9.4.jar
usr/share/kafka/libs/jsr305-3.0.2.jar
usr/share/kafka/libs/kafka-clients-3.9.0.jar
usr/share/kafka/libs/kafka-group-coordinator-3.9.0.jar
usr/share/kafka/libs/kafka-group-coordinator-api-3.9.0.jar
usr/share/kafka/libs/kafka-metadata-3.9.0.jar
usr/share/kafka/libs/kafka-raft-3.9.0.jar
usr/share/kafka/libs/kafka-server-3.9.0.jar
usr/share/kafka/libs/kafka-server-common-3.9.0.jar
usr/share/kafka/libs/kafka-shell-3.9.0.jar
usr/share/kafka/libs/kafka-storage-3.9.0.jar
usr/share/kafka/libs/kafka-storage-api-3.9.0.jar
usr/share/kafka/libs/kafka-streams-3.9.0.jar
usr/share/kafka/libs/kafka-streams-examples-3.9.0.jar
usr/share/kafka/libs/kafka-streams-scala_2.13-3.9.0.jar
usr/share/kafka/libs/kafka-streams-test-utils-3.9.0.jar
usr/share/kafka/libs/kafka-tools-3.9.0.jar
usr/share/kafka/libs/kafka-tools-api-3.9.0.jar
usr/share/kafka/libs/kafka-transaction-coordinator-3.9.0.jar
usr/share/kafka/libs/kafka_2.13-3.9.0.jar
usr/share/kafka/libs/lz4-java-1.8.0.jar
usr/share/kafka/libs/maven-artifact-3.9.6.jar
usr/share/kafka/libs/metrics-core-2.2.0.jar
usr/share/kafka/libs/metrics-core-4.1.12.1.jar
usr/share/kafka/libs/netty-buffer-4.1.111.Final.jar
usr/share/kafka/libs/netty-codec-4.1.111.Final.jar
usr/share/kafka/libs/netty-common-4.1.111.Final.jar
usr/share/kafka/libs/netty-handler-4.1.111.Final.jar
usr/share/kafka/libs/netty-resolver-4.1.111.Final.jar
usr/share/kafka/libs/netty-transport-4.1.111.Final.jar
usr/share/kafka/libs/netty-transport-classes-epoll-4.1.111.Final.jar
usr/share/kafka/libs/netty-transport-native-epoll-4.1.111.Final.jar
usr/share/kafka/libs/netty-transport-native-unix-common-4.1.111.Final.jar
usr/share/kafka/libs/opentelemetry-proto-1.0.0-alpha.jar
usr/share/kafka/libs/osgi-resource-locator-1.0.3.jar
usr/share/kafka/libs/paranamer-2.8.jar
usr/share/kafka/libs/pcollections-4.0.1.jar
usr/share/kafka/libs/plexus-utils-3.5.1.jar
usr/share/kafka/libs/protobuf-java-3.25.5.jar
usr/share/kafka/libs/reflections-0.10.2.jar
usr/share/kafka/libs/reload4j-1.2.25.jar
usr/share/kafka/libs/rocksdbjni-7.9.2.jar
usr/share/kafka/libs/scala-collection-compat_2.13-2.10.0.jar
usr/share/kafka/libs/scala-java8-compat_2.13-1.0.2.jar
usr/share/kafka/libs/scala-library-2.13.14.jar
usr/share/kafka/libs/scala-logging_2.13-3.9.5.jar
usr/share/kafka/libs/scala-reflect-2.13.14.jar
usr/share/kafka/libs/slf4j-api-1.7.36.jar
usr/share/kafka/libs/slf4j-reload4j-1.7.36.jar
usr/share/kafka/libs/snappy-java-1.1.10.5.jar
usr/share/kafka/libs/swagger-annotations-2.2.8.jar
usr/share/kafka/libs/trogdor-3.9.0.jar
usr/share/kafka/libs/zookeeper-3.8.4.jar
usr/share/kafka/libs/zookeeper-jute-3.8.4.jar
usr/share/kafka/libs/zstd-jni-1.5.6-4.jar
usr/share/kafka/licenses/
usr/share/kafka/licenses/CDDL+GPL-1.1
usr/share/kafka/licenses/DWTFYWTPL
usr/share/kafka/licenses/argparse-MIT
usr/share/kafka/licenses/eclipse-distribution-license-1.0
usr/share/kafka/licenses/eclipse-public-license-2.0
usr/share/kafka/licenses/jline-BSD-3-clause
usr/share/kafka/licenses/jopt-simple-MIT
usr/share/kafka/licenses/jsr305-BSD-3-clause
usr/share/kafka/licenses/paranamer-BSD-3-clause
usr/share/kafka/licenses/pcollections-MIT
usr/share/kafka/licenses/protobuf-java-BSD-3-clause
usr/share/kafka/licenses/slf4j-MIT
usr/share/kafka/licenses/zstd-jni-BSD-2-clause
usr/share/kafka/site-docs/
usr/share/kafka/site-docs/kafka_2.13-3.9.0-site-docs.tgz
var/
var/log/
var/log/kafka/
