# extracted/valkey-8.0.1-2/files
%FILES%
etc/
etc/valkey/
etc/valkey/sentinel.conf
etc/valkey/valkey.conf
usr/
usr/bin/
usr/bin/redis-benchmark
usr/bin/redis-check-aof
usr/bin/redis-check-rdb
usr/bin/redis-cli
usr/bin/redis-sentinel
usr/bin/redis-server
usr/bin/valkey-benchmark
usr/bin/valkey-check-aof
usr/bin/valkey-check-rdb
usr/bin/valkey-cli
usr/bin/valkey-sentinel
usr/bin/valkey-server
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/valkey-sentinel.service
usr/lib/systemd/system/valkey.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/valkey.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/valkey.conf
usr/share/
usr/share/licenses/
usr/share/licenses/valkey/
usr/share/licenses/valkey/LICENSE
