# extracted/tinyproxy-1.11.2-2/files
%FILES%
etc/
etc/tinyproxy/
etc/tinyproxy/tinyproxy.conf
usr/
usr/bin/
usr/bin/tinyproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tinyproxy.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/tinyproxy.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/tinyproxy.conf
usr/share/
usr/share/doc/
usr/share/doc/tinyproxy/
usr/share/doc/tinyproxy/AUTHORS
usr/share/doc/tinyproxy/NEWS
usr/share/doc/tinyproxy/README
usr/share/doc/tinyproxy/README.md
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/tinyproxy.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/tinyproxy.8.gz
usr/share/tinyproxy/
usr/share/tinyproxy/debug.html
usr/share/tinyproxy/default.html
usr/share/tinyproxy/stats.html
