# extracted/shairport-sync-4.3.5-1/files
%FILES%
etc/
etc/shairport-sync.conf
usr/
usr/bin/
usr/bin/shairport-sync
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/shairport-sync.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/shairport-sync.conf
usr/share/
usr/share/licenses/
usr/share/licenses/shairport-sync/
usr/share/licenses/shairport-sync/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/shairport-sync.1.gz
