# extracted/dopewars-1.6.2-3/files
%FILES%
srv/
srv/dopewars/
srv/dopewars/dopewars.sco
usr/
usr/bin/
usr/bin/dopewars
usr/lib/
usr/lib/dopewars/
usr/lib/dopewars/libsound_sdl.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dopewars.service
usr/lib/systemd/system/dopewars.socket
usr/share/
usr/share/applications/
usr/share/applications/dopewars.desktop
usr/share/doc/
usr/share/doc/dopewars/
usr/share/doc/dopewars/LICENCE
usr/share/doc/dopewars/aiplayer.html
usr/share/doc/dopewars/clientplay.html
usr/share/doc/dopewars/commandline.html
usr/share/doc/dopewars/configfile.html
usr/share/doc/dopewars/contribute.html
usr/share/doc/dopewars/credits.html
usr/share/doc/dopewars/developer.html
usr/share/doc/dopewars/example-cfg
usr/share/doc/dopewars/example-igneous
usr/share/doc/dopewars/help/
usr/share/doc/dopewars/help/cops.html
usr/share/doc/dopewars/help/drugs.html
usr/share/doc/dopewars/help/general.html
usr/share/doc/dopewars/help/guns.html
usr/share/doc/dopewars/help/locations.html
usr/share/doc/dopewars/help/server.html
usr/share/doc/dopewars/help/sounds.html
usr/share/doc/dopewars/i18n.html
usr/share/doc/dopewars/index.html
usr/share/doc/dopewars/installation.html
usr/share/doc/dopewars/metaserver.html
usr/share/doc/dopewars/protocol.html
usr/share/doc/dopewars/server.html
usr/share/doc/dopewars/servercommands.html
usr/share/doc/dopewars/windows.html
usr/share/dopewars/
usr/share/dopewars/bye.wav
usr/share/dopewars/colt.wav
usr/share/dopewars/die.wav
usr/share/dopewars/gun.wav
usr/share/dopewars/jet.wav
usr/share/dopewars/losebitch.wav
usr/share/dopewars/message.wav
usr/share/dopewars/murmur.wav
usr/share/dopewars/punch.wav
usr/share/dopewars/run.wav
usr/share/dopewars/shotdown.wav
usr/share/dopewars/train.wav
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/dopewars.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/dopewars.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/dopewars.mo
usr/share/locale/es_ES/
usr/share/locale/es_ES/LC_MESSAGES/
usr/share/locale/es_ES/LC_MESSAGES/dopewars.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/dopewars.mo
usr/share/locale/fr_CA/
usr/share/locale/fr_CA/LC_MESSAGES/
usr/share/locale/fr_CA/LC_MESSAGES/dopewars.mo
usr/share/locale/nn/
usr/share/locale/nn/LC_MESSAGES/
usr/share/locale/nn/LC_MESSAGES/dopewars.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/dopewars.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/dopewars.mo
usr/share/man/
usr/share/man/man6/
usr/share/man/man6/dopewars.6.gz
usr/share/pixmaps/
usr/share/pixmaps/dopewars-pill.png
usr/share/pixmaps/dopewars-shot.png
usr/share/pixmaps/dopewars-weed.png
