# extracted/tang-15-1/files
%FILES%
usr/
usr/bin/
usr/bin/tang-show-keys
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tangd.socket
usr/lib/systemd/system/tangd@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/tang.conf
usr/lib/tangd
usr/lib/tangd-keygen
usr/lib/tangd-rotate-keys
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/tang.conf
usr/share/
usr/share/licenses/
usr/share/licenses/tang/
usr/share/licenses/tang/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/tang-show-keys.1.gz
usr/share/man/man1/tangd-rotate-keys.1.gz
usr/share/man/man8/
usr/share/man/man8/tang.8.gz
var/
var/db/
var/db/tang/
