# extracted/igmpproxy-0.4-2/files
%FILES%
etc/
etc/igmpproxy.conf
usr/
usr/bin/
usr/bin/igmpproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/igmpproxy.service
usr/share/
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/igmpproxy.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/igmpproxy.8.gz
