# extracted/zabbix-agent2-7.2.0-1/files
%FILES%
etc/
etc/zabbix/
etc/zabbix/zabbix_agent2.conf
etc/zabbix/zabbix_agent2.d/
etc/zabbix/zabbix_agent2.d/plugins.d/
usr/
usr/bin/
usr/bin/zabbix_agent2
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/zabbix-agent2.service
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/zabbix_agent2.8.gz
