# extracted/keepalived-2.3.2-1/files
%FILES%
etc/
etc/keepalived/
etc/keepalived/keepalived.conf
etc/sysconfig/
etc/sysconfig/keepalived
usr/
usr/bin/
usr/bin/genhash
usr/bin/keepalived
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/keepalived.service
usr/share/
usr/share/doc/
usr/share/doc/keepalived/
usr/share/doc/keepalived/README
usr/share/keepalived/
usr/share/keepalived/samples/
usr/share/keepalived/samples/keepalived.conf.HTTP_GET.port
usr/share/keepalived/samples/keepalived.conf.IPv6
usr/share/keepalived/samples/keepalived.conf.PING_CHECK
usr/share/keepalived/samples/keepalived.conf.SMTP_CHECK
usr/share/keepalived/samples/keepalived.conf.SSL_GET
usr/share/keepalived/samples/keepalived.conf.UDP_CHECK
usr/share/keepalived/samples/keepalived.conf.conditional_conf
usr/share/keepalived/samples/keepalived.conf.fwmark
usr/share/keepalived/samples/keepalived.conf.inhibit
usr/share/keepalived/samples/keepalived.conf.misc_check
usr/share/keepalived/samples/keepalived.conf.misc_check_arg
usr/share/keepalived/samples/keepalived.conf.quorum
usr/share/keepalived/samples/keepalived.conf.sample
usr/share/keepalived/samples/keepalived.conf.status_code
usr/share/keepalived/samples/keepalived.conf.track_interface
usr/share/keepalived/samples/keepalived.conf.virtual_server_group
usr/share/keepalived/samples/keepalived.conf.virtualhost
usr/share/keepalived/samples/keepalived.conf.vrrp
usr/share/keepalived/samples/keepalived.conf.vrrp.localcheck
usr/share/keepalived/samples/keepalived.conf.vrrp.lvs_syncd
usr/share/keepalived/samples/keepalived.conf.vrrp.routes
usr/share/keepalived/samples/keepalived.conf.vrrp.rules
usr/share/keepalived/samples/keepalived.conf.vrrp.scripts
usr/share/keepalived/samples/keepalived.conf.vrrp.static_ipaddress
usr/share/keepalived/samples/keepalived.conf.vrrp.sync
usr/share/keepalived/samples/sample.misccheck.smbcheck.sh
usr/share/keepalived/samples/sample_notify_fifo.sh
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/genhash.1.gz
usr/share/man/man5/
usr/share/man/man5/keepalived.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/keepalived.8.gz
