# extracted/libvirt-dbus-1.4.1-3/files
%FILES%
usr/
usr/bin/
usr/bin/libvirt-dbus
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/libvirt-dbus.service
usr/lib/systemd/user/
usr/lib/systemd/user/libvirt-dbus.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/libvirtdbus.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.libvirt.Connect.xml
usr/share/dbus-1/interfaces/org.libvirt.Domain.xml
usr/share/dbus-1/interfaces/org.libvirt.DomainSnapshot.xml
usr/share/dbus-1/interfaces/org.libvirt.Interface.xml
usr/share/dbus-1/interfaces/org.libvirt.NWFilter.xml
usr/share/dbus-1/interfaces/org.libvirt.Network.xml
usr/share/dbus-1/interfaces/org.libvirt.NodeDevice.xml
usr/share/dbus-1/interfaces/org.libvirt.Secret.xml
usr/share/dbus-1/interfaces/org.libvirt.StoragePool.xml
usr/share/dbus-1/interfaces/org.libvirt.StorageVol.xml
usr/share/dbus-1/services/
usr/share/dbus-1/services/org.libvirt.service
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.libvirt.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.libvirt.conf
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/libvirt-dbus.8.gz
usr/share/polkit-1/
usr/share/polkit-1/rules.d/
usr/share/polkit-1/rules.d/libvirt-dbus.rules
