# extracted/vector-0.43.1-1/files
%FILES%
etc/
etc/default/
etc/default/vector
etc/vector/
etc/vector/vector.yaml
usr/
usr/bin/
usr/bin/vector
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hardened-vector.service
usr/lib/systemd/system/vector.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/vector.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/vector.conf
usr/share/
usr/share/doc/
usr/share/doc/vector/
usr/share/doc/vector/CHANGELOG.md
usr/share/doc/vector/PRIVACY.md
usr/share/doc/vector/README.md
usr/share/doc/vector/RELEASES.md
usr/share/doc/vector/SECURITY.md
usr/share/doc/vector/VERSIONING.md
usr/share/doc/vector/examples/
usr/share/doc/vector/examples/docs_example.yaml
usr/share/doc/vector/examples/environment_variables.yaml
usr/share/doc/vector/examples/es_s3_hybrid.yaml
usr/share/doc/vector/examples/file_to_cloudwatch_metrics.yaml
usr/share/doc/vector/examples/file_to_prometheus.yaml
usr/share/doc/vector/examples/namespacing/
usr/share/doc/vector/examples/namespacing/sinks/
usr/share/doc/vector/examples/namespacing/sinks/es_cluster.yaml
usr/share/doc/vector/examples/namespacing/sinks/s3_archives.yaml
usr/share/doc/vector/examples/namespacing/sources/
usr/share/doc/vector/examples/namespacing/sources/apache_logs.yaml
usr/share/doc/vector/examples/namespacing/transforms/
usr/share/doc/vector/examples/namespacing/transforms/apache_parser.yaml
usr/share/doc/vector/examples/namespacing/transforms/apache_sample.yaml
usr/share/doc/vector/examples/namespacing/vector.yaml
usr/share/doc/vector/examples/prometheus_to_console.yaml
usr/share/doc/vector/examples/stdio.yaml
usr/share/doc/vector/examples/wrapped_json.yaml
usr/share/licenses/
usr/share/licenses/vector/
usr/share/licenses/vector/LICENSE
usr/share/licenses/vector/NOTICE
