# extracted/tailscale-1.78.1-1/files
%FILES%
etc/
etc/default/
etc/default/tailscaled
usr/
usr/bin/
usr/bin/tailscale
usr/bin/tailscaled
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tailscaled.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/tailscale
usr/share/fish/
usr/share/fish/vendor_completions.d/
usr/share/fish/vendor_completions.d/tailscale.fish
usr/share/licenses/
usr/share/licenses/tailscale/
usr/share/licenses/tailscale/LICENSE
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_tailscale
