# extracted/ucarp-1.5.2-8/files
%FILES%
usr/
usr/bin/
usr/bin/ucarp
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ucarp.service
usr/share/
usr/share/doc/
usr/share/doc/ucarp/
usr/share/doc/ucarp/README
usr/share/doc/ucarp/examples/
usr/share/doc/ucarp/examples/vip-down.sh
usr/share/doc/ucarp/examples/vip-up.sh
usr/share/locale/
usr/share/locale/en@boldquot/
usr/share/locale/en@boldquot/LC_MESSAGES/
usr/share/locale/en@boldquot/LC_MESSAGES/ucarp.mo
usr/share/locale/en@quot/
usr/share/locale/en@quot/LC_MESSAGES/
usr/share/locale/en@quot/LC_MESSAGES/ucarp.mo
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/ucarp.8.gz
