# extracted/opensnitch-1.6.6-3/files
%FILES%
etc/
etc/logrotate.d/
etc/logrotate.d/opensnitch
etc/opensnitchd/
etc/opensnitchd/default-config.json
etc/opensnitchd/system-fw.json
usr/
usr/bin/
usr/bin/opensnitch-ui
usr/bin/opensnitchd
usr/lib/
usr/lib/python3.13/
usr/lib/python3.13/site-packages/
usr/lib/python3.13/site-packages/opensnitch/
usr/lib/python3.13/site-packages/opensnitch/__init__.py
usr/lib/python3.13/site-packages/opensnitch/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/config.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/config.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/desktop_parser.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/desktop_parser.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/nodes.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/nodes.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/notifications.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/notifications.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/resources_rc.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/resources_rc.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/rules.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/rules.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/service.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/service.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/ui_pb2.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/ui_pb2.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/ui_pb2_grpc.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/ui_pb2_grpc.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/version.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/__pycache__/version.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/
usr/lib/python3.13/site-packages/opensnitch/actions/__init__.py
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/default_configs.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/default_configs.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/highlight.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/highlight.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/utils.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/__pycache__/utils.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/actions/default_configs.py
usr/lib/python3.13/site-packages/opensnitch/actions/highlight.py
usr/lib/python3.13/site-packages/opensnitch/actions/utils.py
usr/lib/python3.13/site-packages/opensnitch/auth/
usr/lib/python3.13/site-packages/opensnitch/auth/__init__.py
usr/lib/python3.13/site-packages/opensnitch/auth/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/auth/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/auth/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/config.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__init__.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/addresstablemodel.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/addresstablemodel.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/colorizeddelegate.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/colorizeddelegate.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/firewalltableview.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/firewalltableview.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/generictableview.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/generictableview.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/main.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/main.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/updownbtndelegate.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/__pycache__/updownbtndelegate.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/customwidgets/addresstablemodel.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/colorizeddelegate.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/firewalltableview.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/generictableview.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/main.py
usr/lib/python3.13/site-packages/opensnitch/customwidgets/updownbtndelegate.py
usr/lib/python3.13/site-packages/opensnitch/database/
usr/lib/python3.13/site-packages/opensnitch/database/__init__.py
usr/lib/python3.13/site-packages/opensnitch/database/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/database/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/database/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/database/__pycache__/enums.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/database/__pycache__/enums.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/database/enums.py
usr/lib/python3.13/site-packages/opensnitch/database/migrations/
usr/lib/python3.13/site-packages/opensnitch/database/migrations/upgrade_1.sql
usr/lib/python3.13/site-packages/opensnitch/database/migrations/upgrade_2.sql
usr/lib/python3.13/site-packages/opensnitch/database/migrations/upgrade_3.sql
usr/lib/python3.13/site-packages/opensnitch/desktop_parser.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/
usr/lib/python3.13/site-packages/opensnitch/dialogs/__init__.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/conndetails.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/conndetails.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/firewall.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/firewall.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/firewall_rule.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/firewall_rule.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/preferences.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/preferences.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/processdetails.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/processdetails.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/prompt.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/prompt.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/ruleseditor.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/ruleseditor.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/stats.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/__pycache__/stats.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/dialogs/conndetails.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/firewall.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/firewall_rule.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/preferences.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/processdetails.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/prompt.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/ruleseditor.py
usr/lib/python3.13/site-packages/opensnitch/dialogs/stats.py
usr/lib/python3.13/site-packages/opensnitch/firewall/
usr/lib/python3.13/site-packages/opensnitch/firewall/__init__.py
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/chains.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/chains.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/enums.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/enums.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/exprs.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/exprs.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/profiles.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/profiles.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/rules.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/rules.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/utils.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/__pycache__/utils.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/firewall/chains.py
usr/lib/python3.13/site-packages/opensnitch/firewall/enums.py
usr/lib/python3.13/site-packages/opensnitch/firewall/exprs.py
usr/lib/python3.13/site-packages/opensnitch/firewall/profiles.py
usr/lib/python3.13/site-packages/opensnitch/firewall/rules.py
usr/lib/python3.13/site-packages/opensnitch/firewall/utils.py
usr/lib/python3.13/site-packages/opensnitch/i18n/
usr/lib/python3.13/site-packages/opensnitch/i18n/de_DE/
usr/lib/python3.13/site-packages/opensnitch/i18n/de_DE/opensnitch-de_DE.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/es_ES/
usr/lib/python3.13/site-packages/opensnitch/i18n/es_ES/opensnitch-es_ES.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/eu_ES/
usr/lib/python3.13/site-packages/opensnitch/i18n/eu_ES/opensnitch-eu_ES.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/fi_FI/
usr/lib/python3.13/site-packages/opensnitch/i18n/fi_FI/opensnitch-fi_FI.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/fr_FR/
usr/lib/python3.13/site-packages/opensnitch/i18n/fr_FR/opensnitch-fr_FR.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/hu_HU/
usr/lib/python3.13/site-packages/opensnitch/i18n/hu_HU/opensnitch-hu_HU.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/ja_JP/
usr/lib/python3.13/site-packages/opensnitch/i18n/ja_JP/opensnitch-ja_JP.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/lt_LT/
usr/lib/python3.13/site-packages/opensnitch/i18n/lt_LT/opensnitch-lt_LT.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/nb_NO/
usr/lib/python3.13/site-packages/opensnitch/i18n/nb_NO/opensnitch-nb_NO.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/nl_NL/
usr/lib/python3.13/site-packages/opensnitch/i18n/nl_NL/opensnitch-nl_NL.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/pt_BR/
usr/lib/python3.13/site-packages/opensnitch/i18n/pt_BR/opensnitch-pt_BR.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/ro_RO/
usr/lib/python3.13/site-packages/opensnitch/i18n/ro_RO/opensnitch-ro_RO.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/ru_RU/
usr/lib/python3.13/site-packages/opensnitch/i18n/ru_RU/opensnitch-ru_RU.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/tr_TR/
usr/lib/python3.13/site-packages/opensnitch/i18n/tr_TR/opensnitch-tr_TR.qm
usr/lib/python3.13/site-packages/opensnitch/i18n/zh_TW/
usr/lib/python3.13/site-packages/opensnitch/i18n/zh_TW/opensnitch-zh_TW.qm
usr/lib/python3.13/site-packages/opensnitch/nodes.py
usr/lib/python3.13/site-packages/opensnitch/notifications.py
usr/lib/python3.13/site-packages/opensnitch/res/
usr/lib/python3.13/site-packages/opensnitch/res/__init__.py
usr/lib/python3.13/site-packages/opensnitch/res/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/res/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/res/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/res/firewall.ui
usr/lib/python3.13/site-packages/opensnitch/res/firewall_rule.ui
usr/lib/python3.13/site-packages/opensnitch/res/icon-alert.png
usr/lib/python3.13/site-packages/opensnitch/res/icon-off.png
usr/lib/python3.13/site-packages/opensnitch/res/icon-pause.png
usr/lib/python3.13/site-packages/opensnitch/res/icon-pause.svg
usr/lib/python3.13/site-packages/opensnitch/res/icon-red.png
usr/lib/python3.13/site-packages/opensnitch/res/icon-white.png
usr/lib/python3.13/site-packages/opensnitch/res/icon-white.svg
usr/lib/python3.13/site-packages/opensnitch/res/icon.png
usr/lib/python3.13/site-packages/opensnitch/res/preferences.ui
usr/lib/python3.13/site-packages/opensnitch/res/process_details.ui
usr/lib/python3.13/site-packages/opensnitch/res/prompt.ui
usr/lib/python3.13/site-packages/opensnitch/res/resources.qrc
usr/lib/python3.13/site-packages/opensnitch/res/ruleseditor.ui
usr/lib/python3.13/site-packages/opensnitch/res/stats.ui
usr/lib/python3.13/site-packages/opensnitch/resources_rc.py
usr/lib/python3.13/site-packages/opensnitch/rules.py
usr/lib/python3.13/site-packages/opensnitch/service.py
usr/lib/python3.13/site-packages/opensnitch/ui_pb2.py
usr/lib/python3.13/site-packages/opensnitch/ui_pb2_grpc.py
usr/lib/python3.13/site-packages/opensnitch/utils/
usr/lib/python3.13/site-packages/opensnitch/utils/__init__.py
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/__init__.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/__init__.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/infowindow.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/infowindow.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/languages.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/languages.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/xdg.cpython-313.opt-1.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/__pycache__/xdg.cpython-313.pyc
usr/lib/python3.13/site-packages/opensnitch/utils/infowindow.py
usr/lib/python3.13/site-packages/opensnitch/utils/languages.py
usr/lib/python3.13/site-packages/opensnitch/utils/xdg.py
usr/lib/python3.13/site-packages/opensnitch/version.py
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/LICENSE
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/METADATA
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/RECORD
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/WHEEL
usr/lib/python3.13/site-packages/opensnitch_ui-1.6.6.dist-info/top_level.txt
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/opensnitchd.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/opensnitch.conf
usr/share/
usr/share/applications/
usr/share/applications/opensnitch_ui.desktop
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/scalable/
usr/share/icons/hicolor/scalable/apps/
usr/share/icons/hicolor/scalable/apps/opensnitch-ui.svg
usr/share/kservices5/
usr/share/kservices5/kcm_opensnitch.desktop
usr/share/metainfo/
usr/share/metainfo/io.github.evilsocket.opensnitch.appdata.xml
