# extracted/tenshi-0.18-1/files
%FILES%
usr/
usr/bin/
usr/bin/tenshi
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tenshi.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/tenshi.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/tenshi.conf
usr/share/
usr/share/doc/
usr/share/doc/tenshi/
usr/share/doc/tenshi/Changelog
usr/share/doc/tenshi/README.md
usr/share/factory/
usr/share/factory/etc/
usr/share/factory/etc/tenshi/
usr/share/factory/etc/tenshi/tenshi.conf
usr/share/licenses/
usr/share/licenses/tenshi/
usr/share/licenses/tenshi/LICENSE
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/tenshi.8.gz
