# extracted/prometheus-postgres-exporter-0.16.0-1/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/prometheus-postgres-exporter
usr/
usr/bin/
usr/bin/prometheus-postgres-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-postgres-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-postgres-exporter.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/prometheus-postgres-exporter.conf
usr/share/
usr/share/doc/
usr/share/doc/prometheus-postgres-exporter/
usr/share/doc/prometheus-postgres-exporter/README.md
usr/share/licenses/
usr/share/licenses/prometheus-postgres-exporter/
usr/share/licenses/prometheus-postgres-exporter/LICENSE
