# extracted/isatapd-0.9.7-8/files
%FILES%
usr/
usr/bin/
usr/bin/isatapd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/isatapd@.service
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/isatapd.8.gz
