# extracted/rng-tools-6.17-2/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/rngd
usr/
usr/bin/
usr/bin/randstat
usr/bin/rngd
usr/bin/rngtest
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rngd.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/rngtest.1.gz
usr/share/man/man8/
usr/share/man/man8/rngd.8.gz
