# extracted/chrony-4.6.1-2/files
%FILES%
etc/
etc/chrony.conf
etc/logrotate.d/
etc/logrotate.d/chrony
usr/
usr/bin/
usr/bin/chronyc
usr/bin/chronyd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/ntp-units.d/
usr/lib/systemd/ntp-units.d/50-chrony.list
usr/lib/systemd/system/
usr/lib/systemd/system/chrony-wait.service
usr/lib/systemd/system/chronyd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/chrony.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/chrony.conf
usr/share/
usr/share/doc/
usr/share/doc/chrony/
usr/share/doc/chrony/NEWS
usr/share/doc/chrony/README
usr/share/doc/chrony/chrony.conf.html
usr/share/doc/chrony/chronyc.html
usr/share/doc/chrony/chronyd.html
usr/share/doc/chrony/examples/
usr/share/doc/chrony/examples/chrony-wait.service
usr/share/doc/chrony/examples/chrony.conf.example1
usr/share/doc/chrony/examples/chrony.conf.example2
usr/share/doc/chrony/examples/chrony.conf.example3
usr/share/doc/chrony/examples/chrony.keys.example
usr/share/doc/chrony/examples/chrony.logrotate
usr/share/doc/chrony/examples/chrony.nm-dispatcher.dhcp
usr/share/doc/chrony/examples/chrony.nm-dispatcher.onoffline
usr/share/doc/chrony/examples/chronyd-restricted.service
usr/share/doc/chrony/examples/chronyd.service
usr/share/doc/chrony/faq.html
usr/share/doc/chrony/installation.html
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/chronyc.1.gz
usr/share/man/man5/
usr/share/man/man5/chrony.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/chronyd.8.gz
