# extracted/ptunnel-0.72-6/files
%FILES%
usr/
usr/bin/
usr/bin/ptunnel
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ptunnel.service
usr/share/
usr/share/licenses/
usr/share/licenses/ptunnel/
usr/share/licenses/ptunnel/LICENSE
