# extracted/traefik-3.2.1-1/files
%FILES%
etc/
etc/traefik/
etc/traefik/acme.json
etc/traefik/traefik.yml
usr/
usr/bin/
usr/bin/traefik
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/traefik.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/traefik.conf
usr/share/
usr/share/licenses/
usr/share/licenses/traefik/
usr/share/licenses/traefik/LICENSE.md
