# extracted/prometheus-ipmi-exporter-1.8.0-3/files
%FILES%
etc/
etc/prometheus/
etc/prometheus/ipmi-exporter.env
etc/prometheus/ipmi-exporter.yml
etc/sudoers.d/
etc/sudoers.d/prometheus-ipmi-exporter
usr/
usr/bin/
usr/bin/prometheus-ipmi-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-ipmi-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-ipmi-exporter.conf
usr/share/
usr/share/doc/
usr/share/doc/prometheus-ipmi-exporter/
usr/share/doc/prometheus-ipmi-exporter/README.md
usr/share/doc/prometheus-ipmi-exporter/configuration.md
usr/share/doc/prometheus-ipmi-exporter/ipmi_local.yml
usr/share/doc/prometheus-ipmi-exporter/ipmi_remote.yml
usr/share/doc/prometheus-ipmi-exporter/metrics.md
usr/share/licenses/
usr/share/licenses/prometheus-ipmi-exporter/
usr/share/licenses/prometheus-ipmi-exporter/LICENSE
