# extracted/prometheus-ssl-exporter-2.4.3-2/files
%FILES%
etc/
etc/prometheus/
etc/prometheus/ssl-exporter.env
etc/prometheus/ssl-exporter.yml
usr/
usr/bin/
usr/bin/prometheus-ssl-exporter
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/prometheus-ssl-exporter.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/prometheus-ssl-exporter.conf
usr/share/
usr/share/licenses/
usr/share/licenses/prometheus-ssl-exporter/
usr/share/licenses/prometheus-ssl-exporter/LICENSE
