# extracted/luanti-server-5.10.0-1/files
%FILES%
etc/
etc/luanti/
usr/
usr/bin/
usr/bin/luantiserver
usr/bin/minetestserver
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/luanti@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/luanti-server.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/luanti-server.conf
usr/share/
usr/share/man/
usr/share/man/man6/
usr/share/man/man6/luantiserver.6.gz
