# extracted/transmission-cli-4.0.6-2/files
%FILES%
usr/
usr/bin/
usr/bin/transmission-cli
usr/bin/transmission-create
usr/bin/transmission-daemon
usr/bin/transmission-edit
usr/bin/transmission-remote
usr/bin/transmission-show
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/transmission.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/transmission.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/transmission.conf
usr/share/
usr/share/licenses/
usr/share/licenses/transmission-cli/
usr/share/licenses/transmission-cli/COPYING
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/transmission-cli.1.gz
usr/share/man/man1/transmission-create.1.gz
usr/share/man/man1/transmission-daemon.1.gz
usr/share/man/man1/transmission-edit.1.gz
usr/share/man/man1/transmission-remote.1.gz
usr/share/man/man1/transmission-show.1.gz
usr/share/transmission/
usr/share/transmission/public_html/
usr/share/transmission/public_html/images/
usr/share/transmission/public_html/images/favicon.ico
usr/share/transmission/public_html/images/favicon.png
usr/share/transmission/public_html/images/webclip-icon.png
usr/share/transmission/public_html/index.html
usr/share/transmission/public_html/transmission-app.css
usr/share/transmission/public_html/transmission-app.js
usr/share/transmission/public_html/transmission-app.js.LEGAL.txt
usr/share/transmission/public_html/transmission-app.js.LICENSE.txt
usr/share/transmission/public_html/transmission-app.js.map
