# extracted/wireproxy-1.0.9-1/files
%FILES%
usr/
usr/bin/
usr/bin/wireproxy
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wireproxy.service
usr/lib/systemd/system/wireproxy@.service
usr/lib/systemd/user/
usr/lib/systemd/user/wireproxy.service
usr/lib/systemd/user/wireproxy@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/wireproxy.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/wireproxy.conf
usr/share/
usr/share/licenses/
usr/share/licenses/wireproxy/
usr/share/licenses/wireproxy/LICENSE
