# extracted/rebuilderd-0.21.0-1/files
%FILES%
etc/
etc/rebuilderd-sync.conf
etc/rebuilderd-worker.conf
etc/rebuilderd.conf
usr/
usr/bin/
usr/bin/rebuilderd
usr/bin/rebuilderd-worker
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rebuilderd-sync@.service
usr/lib/systemd/system/rebuilderd-sync@.timer
usr/lib/systemd/system/rebuilderd-worker@.service
usr/lib/systemd/system/rebuilderd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/rebuilderd.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rebuilderd.conf
usr/libexec/
usr/libexec/rebuilderd/
usr/libexec/rebuilderd/rebuilder-archlinux.sh
usr/libexec/rebuilderd/rebuilder-debian.sh
usr/libexec/rebuilderd/rebuilder-tails.sh
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/rebuilderd-worker.1.gz
usr/share/man/man1/rebuilderd.1.gz
usr/share/man/man5/
usr/share/man/man5/rebuilderd-sync.conf.5.gz
usr/share/man/man5/rebuilderd-worker.conf.5.gz
usr/share/man/man5/rebuilderd.conf.5.gz
