# extracted/rustypaste-0.16.0-1/files
%FILES%
etc/
etc/rustypaste/
etc/rustypaste/config.toml
etc/rustypaste/rustypaste.env
usr/
usr/bin/
usr/bin/rustypaste
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rustypaste.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/rustypaste.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rustypaste.conf
usr/share/
usr/share/doc/
usr/share/doc/rustypaste/
usr/share/doc/rustypaste/README.md
usr/share/licenses/
usr/share/licenses/rustypaste/
usr/share/licenses/rustypaste/LICENSE
